# intel/.llm_actor_uc_cache/23/2366f9eab0f8999388a1c01dee0f144173a282cb.json
{
  "name": "BITTER",
  "ucs": [
    {
      "title": "[LLM] BITTER (T-APT-17) Equation Editor exploit chain \u2014 EQNEDT32.EXE spawns shell/downloader (CVE-2017-11882 / CVE-2018-0798 weaponized RTF)",
      "description": "Hunts BITTER's signature initial-access chain against South-Asian government/energy targets: a weaponized RTF (delivered via spearphish to Pakistan/Bangladesh/China defense and engineering verticals) abuses the unpatched Microsoft Equation Editor (EQNEDT32.EXE) to execute shellcode that drops ArtraDownloader / ZxxZ / MuuyDownloader / WmRAT. EQNEDT32 is OOP-invoked by Office and has no legitimate reason to spawn a child \u2014 *any* child is a high-fidelity hit on this group's tradecraft (Forcepoint 2016, Talos 2022 Bangladesh campaign, Proofpoint 2024 WmRAT/MiyaRAT cluster).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "exploit",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1203",
        "T1204.002",
        "T1059.003"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.parent_process_name=\"EQNEDT32.EXE\" by host Processes.dest Processes.user Processes.parent_process_name Processes.parent_process Processes.process_name Processes.process Processes.process_hash | `drop_dm_object_name(Processes)` | rename dest as Computer | where NOT match(user, \"\\$$\") | eval risk_score=90, threat_actor=\"BITTER (G1002 / T-APT-17)\", cve=\"CVE-2017-11882, CVE-2018-0798\" | convert ctime(firstTime) ctime(lastTime) | sort - lastTime",
      "kql": "// BITTER (G1002) \u2014 Equation Editor exploitation chain\n// EQNEDT32.EXE has no legitimate reason to spawn child processes\nlet _suspicious_children = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"curl.exe\",\"wget.exe\",\"schtasks.exe\",\"reg.exe\",\"forfiles.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(30d)\n| where InitiatingProcessFileName =~ \"eqnedt32.exe\"\n| where AccountName !endswith \"$\"\n| extend SuspiciousChild = FileName in~ (_suspicious_children)\n| project Timestamp, DeviceName, AccountName,\n          ParentPath = InitiatingProcessFolderPath,\n          ParentCmd  = InitiatingProcessCommandLine,\n          ParentSHA  = InitiatingProcessSHA256,\n          Child      = FileName,\n          ChildPath  = FolderPath,\n          ChildCmd   = ProcessCommandLine,\n          ChildSHA   = SHA256,\n          SuspiciousChild,\n          IsRemoteSession = InitiatingProcessTokenElevation\n| order by Timestamp desc",
      "sentinel_kql": "// BITTER (G1002) \u2014 Equation Editor exploitation chain (Sentinel SecurityEvent EID 4688)\nlet _suspicious_children = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"curl.exe\",\"wget.exe\",\"schtasks.exe\",\"reg.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(30d)\n| where EventID == 4688\n| where ParentProcessName endswith \"EQNEDT32.EXE\"\n| where Account !endswith \"$\"\n| extend ChildBin = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n| extend SuspiciousChild = ChildBin in~ (_suspicious_children)\n| project TimeGenerated, Computer, Account,\n          ParentProcessName, NewProcessName, CommandLine,\n          TokenElevationType, SuspiciousChild\n| order by TimeGenerated desc",
      "sigma_yaml": "title: BITTER APT - Microsoft Equation Editor Spawns Child Process\nid: 9c4ad2e1-6f3b-4a7d-8e5f-1b2c3d4e5f6a\nstatus: experimental\ndescription: |\n  Detects EQNEDT32.EXE spawning any child process. Microsoft Equation Editor was retired\n  in 2018 and is exclusively reached via embedded OLE in malicious RTF/DOC files. BITTER\n  (G1002 / T-APT-17) operationalises CVE-2017-11882 and CVE-2018-0798 against Pakistan,\n  Bangladesh, and Chinese government/energy targets to drop ArtraDownloader, ZxxZ,\n  MuuyDownloader, and WmRAT.\nreferences:\n  - https://attack.mitre.org/groups/G1002/\n  - https://blog.talosintelligence.com/2022/05/bitter-apt-adds-bangladesh-to-their.html\n  - https://www.forcepoint.com/blog/x-labs/bitter-targeted-attack-against-pakistan\nauthor: SOC Detection Engineering\ndate: 2026/05/05\ntags:\n  - attack.initial_access\n  - attack.execution\n  - attack.t1566.001\n  - attack.t1203\n  - attack.t1204.002\n  - attack.g1002\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection:\n    ParentImage|endswith: '\\EQNEDT32.EXE'\n  condition: selection\nfalsepositives:\n  - Equation Editor support ended in January 2018; legitimate launches in modern enterprises are essentially zero.\nlevel: high",
      "rationale": "BITTER has consistently weaponized Equation Editor RCEs (CVE-2017-11882 confirmed in the Talos 2022 Bangladesh military report; CVE-2018-0798 used in 2020-2024 RTF lures). EQNEDT32.EXE is an out-of-process COM server with no scripting/automation customers \u2014 a child process is essentially pathognomonic of exploitation. Other South-Asian-aligned actors (SideWinder, Patchwork, Donot) also touch these"
    },
    {
      "title": "[LLM] BITTER scheduled-task downloader persistence \u2014 schtasks /create invoking LOLBin to fetch from DDNS C2",
      "description": "Hunts BITTER's distinctive post-exploit persistence pattern (Talos, Proofpoint, Cisco Secureworks): schtasks.exe is used to register a recurring task that fires curl/PowerShell/mshta/certutil against a Dynamic DNS host (No-IP, ChangeIP, DuckDNS) to pull stage-2 (BitterRAT, Almond RAT, WmRAT, MiyaRAT). The combination of schtasks.exe + an inline LOLBin URL fetch + a DDNS provider in the same command line is rarely benign and matches BITTER's kill chain documented across G1002 reporting.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1053.005",
        "T1105",
        "T1568",
        "T1071.001",
        "T1036.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" Processes.process=\"*/create*\" (Processes.process=\"*powershell*\" OR Processes.process=\"*cmd /c*\" OR Processes.process=\"*mshta*\" OR Processes.process=\"*curl*\" OR Processes.process=\"*bitsadmin*\" OR Processes.process=\"*certutil*\" OR Processes.process=\"*rundll32*\" OR Processes.process=\"*regsvr32*\" OR Processes.process=\"*wscript*\" OR Processes.process=\"*cscript*\" OR Processes.process=\"*http://*\" OR Processes.process=\"*https://*\" OR Processes.process=\"*ddns.net*\" OR Processes.process=\"*serveftp.com*\" OR Processes.process=\"*linkpc.net*\" OR Processes.process=\"*hopto.org*\" OR Processes.process=\"*sytes.net*\" OR Processes.process=\"*zapto.org*\" OR Processes.process=\"*myftp.org*\" OR Processes.process=\"*duckdns.org*\" OR Processes.process=\"*redirectme.net*\" OR Processes.process=\"*servehttp.com*\" OR Processes.process=\"*ignorelist.com*\") by host Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process | `drop_dm_object_name(Processes)` | rename dest as Computer | where NOT match(user, \"\\$$\") | eval threat_actor=\"BITTER (G1002 / T-APT-17)\" | convert ctime(firstTime) ctime(lastTime) | sort - lastTime",
      "kql": "// BITTER (G1002) \u2014 schtasks /create + LOLBin URL/DDNS fetch persistence\nlet _bitter_ddns = dynamic([\n    \"ddns.net\",\"serveftp.com\",\"linkpc.net\",\"hopto.org\",\"sytes.net\",\n    \"zapto.org\",\"myftp.org\",\"duckdns.org\",\"redirectme.net\",\"servehttp.com\",\n    \"servepics.com\",\"servegame.com\",\"myvnc.com\",\"ignorelist.com\",\n    \"changeip.com\",\"hopto.me\",\"chickenkiller.com\"\n]);\nlet _lolbin_tokens = dynamic([\n    \"powershell\",\"pwsh\",\"cmd /c\",\"cmd.exe /c\",\"mshta\",\"curl\",\"wget\",\n    \"bitsadmin\",\"certutil\",\"rundll32\",\"regsvr32\",\"wscript\",\"cscript\",\n    \"DownloadString\",\"DownloadFile\",\"Invoke-WebRequest\",\"IEX\",\n    \"http://\",\"https://\"\n]);\nDeviceProcessEvents\n| where Timestamp > ago(30d)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has \"/create\" or ProcessCommandLine has \"/CREATE\"\n| where AccountName !endswith \"$\"\n| extend CmdLower = tolower(ProcessCommandLine)\n| where CmdLower has_any (_lolbin_tokens) or CmdLower has_any (_bitter_ddns)\n| extend HasDDNS = CmdLower has_any (_bitter_ddns)\n| project Timestamp, DeviceName, AccountName,\n          TaskCmd     = ProcessCommandLine,\n          Parent      = InitiatingProcessFileName,\n          ParentCmd   = InitiatingProcessCommandLine,\n          ParentPath  = InitiatingProcessFolderPath,\n          GrandParent = InitiatingProcessParentFileName,\n          HasDDNS,\n          IntegrityLvl = ProcessIntegrityLevel\n| order by Timestamp desc",
      "sentinel_kql": "// BITTER (G1002) \u2014 schtasks /create + LOLBin URL/DDNS fetch (Sentinel SecurityEvent EID 4688)\nlet _bitter_ddns = dynamic([\n    \"ddns.net\",\"serveftp.com\",\"linkpc.net\",\"hopto.org\",\"sytes.net\",\n    \"zapto.org\",\"myftp.org\",\"duckdns.org\",\"redirectme.net\",\"servehttp.com\",\n    \"changeip.com\",\"chickenkiller.com\",\"ignorelist.com\"\n]);\nlet _lolbin_tokens = dynamic([\n    \"powershell\",\"pwsh\",\"cmd /c\",\"mshta\",\"curl\",\"wget\",\"bitsadmin\",\n    \"certutil\",\"rundll32\",\"regsvr32\",\"wscript\",\"cscript\",\n    \"downloadstring\",\"downloadfile\",\"invoke-webrequest\",\"http://\",\"https://\"\n]);\nSecurityEvent\n| where TimeGenerated > ago(30d)\n| where EventID == 4688\n| where NewProcessName endswith \"schtasks.exe\"\n| where Account !endswith \"$\"\n| extend CmdLower = tolower(CommandLine)\n| where CmdLower has \"/create\"\n| where CmdLower has_any (_lolbin_tokens) or CmdLower has_any (_bitter_ddns)\n| extend HasDDNS = CmdLower has_any (_bitter_ddns)\n| project TimeGenerated, Computer, Account,\n          CommandLine, ParentProcessName, NewProcessName,\n          HasDDNS, TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "BITTER's downloader stages (ArtraDownloader, MuuyDownloader, the WmRAT/MiyaRAT installer chain) all converge on a near-identical persistence primitive: a scheduled task whose action is a one-liner LOLBin reaching out to a No-IP/ChangeIP/DuckDNS host. The actor relies on DDNS rather than acquired domains because of operational tempo (T1568 + T1583.001 in their MITRE list) and abuses native schtasks"
    }
  ]
}