# intel/.llm_actor_uc_cache/33/3362f34c04e96cc792a0c009d75f8907ffd36abe.json
{
  "name": "TA2541",
  "ucs": [
    {
      "title": "[LLM] TA2541 aviation-lure \u2192 cloud-hosted VBS/MSHTA stager chain (Proofpoint cluster)",
      "description": "Hunts TA2541's signature delivery chain: aviation/aerospace/cargo-themed phishing links pointing to commodity cloud hosts (Google Drive, OneDrive, Discord CDN, MediaFire, GitHub) that drop an obfuscated VBS or LNK which spawns mshta/wscript/powershell to pull a crypter-wrapped commodity RAT (AsyncRAT, NetWire, WSH RAT, Parallax, Revenge RAT). This combination of cloud-hosted aviation lure + scripting-host download is characteristic of TA2541 / Operation Layover (Talos), distinguishing them from generic commodity-RAT operators.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1566.002",
        "T1204.001",
        "T1059.005",
        "T1059.001",
        "T1218.005",
        "T1105",
        "T1583.006",
        "T1608.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"outlook.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"winrar.exe\",\"7zg.exe\",\"explorer.exe\") AND Processes.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\") AND (Processes.process IN (\"*drive.google.com*\",\"*docs.google.com*\",\"*onedrive.live.com*\",\"*1drv.ms*\",\"*cdn.discordapp.com*\",\"*mediafire.com*\",\"*raw.githubusercontent.com*\",\"*pastebin.com*\",\"*bit.ly*\",\"*hxxp*\") OR Processes.process IN (\"*Invoice*\",\"*Cargo*\",\"*Charter*\",\"*Flight*\",\"*Aviation*\",\"*Boeing*\",\"*Airbus*\",\"*Aircraft*\",\"*Quotation*\",\"*RFQ*\",\"*Manifest*\",\"*BL*.vbs\",\"*.vbs\",\"*.hta\",\"*.lnk\")) by Processes.dest Processes.user Processes.process_name Processes.parent_process_name Processes.process \n| `drop_dm_object_name(Processes)` \n| join type=inner dest [| tstats `summariesonly` count from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest_category=\"ddns\" OR All_Traffic.dest IN (\"*duckdns.org\",\"*hopto.org\",\"*ddns.net\",\"*no-ip.com\",\"*ydns.eu\",\"*linkpc.net\") by All_Traffic.dest All_Traffic.src | rename All_Traffic.src as dest] \n| eval risk_score=80, mitre=\"T1566.002+T1218.005+T1059.005\" \n| table firstTime lastTime dest user parent process_name process risk_score",
      "kql": "let aviationLures = dynamic([\"invoice\",\"cargo\",\"charter\",\"flight\",\"aviation\",\"boeing\",\"airbus\",\"aircraft\",\"quotation\",\"rfq\",\"manifest\",\"airway\",\"shipment\",\"freight\",\"trip\",\"itinerary\"]);\nlet cloudHosts = dynamic([\"drive.google.com\",\"docs.google.com\",\"onedrive.live.com\",\"1drv.ms\",\"cdn.discordapp.com\",\"mediafire.com\",\"raw.githubusercontent.com\",\"pastebin.com\",\"transfer.sh\",\"anonfiles.com\"]);\nlet ddnsHosts = dynamic([\".duckdns.org\",\".hopto.org\",\".ddns.net\",\".no-ip.com\",\".ydns.eu\",\".linkpc.net\",\".sytes.net\",\".zapto.org\"]);\nlet stage1 =\nDeviceProcessEvents\n| where Timestamp > ago(14d)\n| where FileName in~ (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\")\n| where InitiatingProcessFileName in~ (\"outlook.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"winrar.exe\",\"7zg.exe\",\"explorer.exe\")\n| where ProcessCommandLine has_any (cloudHosts) or ProcessCommandLine has_any (aviationLures) or ProcessCommandLine has_any (\".vbs\",\".hta\",\".lnk\",\".js\")\n| project Stage1Time=Timestamp, DeviceId, DeviceName, AccountName, ParentProc=InitiatingProcessFileName, ScriptHost=FileName, Stage1Cmd=ProcessCommandLine;\nlet stage2 =\nDeviceNetworkEvents\n| where Timestamp > ago(14d)\n| where RemoteUrl has_any (ddnsHosts) or RemoteUrl has_any (cloudHosts)\n| project Stage2Time=Timestamp, DeviceId, RemoteUrl, RemoteIP, RemotePort, NetInitiator=InitiatingProcessFileName;\nstage1\n| join kind=inner stage2 on DeviceId\n| where Stage2Time between (Stage1Time .. Stage1Time + 10m)\n| extend EmailLure = tostring(extract(@\"(invoice|cargo|charter|flight|aviation|boeing|airbus|aircraft|quotation|rfq|manifest)\", 1, tolower(Stage1Cmd)))\n| join kind=leftouter (EmailEvents | where Timestamp > ago(14d) | where Subject has_any (aviationLures) | project EmailTime=Timestamp, RecipientEmailAddress, SenderFromAddress, Subject, NetworkMessageId) on $left.AccountName == $right.RecipientEmailAddress\n| project Stage1Time, DeviceName, AccountName, ParentProc, ScriptHost, Stage1Cmd, RemoteUrl, RemoteIP, NetInitiator, SenderFromAddress, Subject\n| order by Stage1Time desc",
      "rationale": "TA2541 (Proofpoint) / Operation Layover (Talos) consistently targets aviation, aerospace, defense, transport and manufacturing with aviation-themed lures (RFQs, charter requests, cargo manifests) linking to cloud-hosted VBS/HTA/LNK loaders that beacon to DDNS providers. The join of (mail-client/browser \u2192 script host \u2192 cloud-host URL) followed within minutes by DDNS egress is much more specific tha"
    },
    {
      "title": "[LLM] TA2541 commodity-RAT staging via process hollowing into .NET host binaries (RegSvcs / AspNetCompiler / InstallUtil)",
      "description": "Detects TA2541's post-exploitation hallmark: a crypter spawned by wscript/powershell/mshta that hollows a benign signed .NET utility (RegSvcs.exe, RegAsm.exe, InstallUtil.exe, AspNetCompiler.exe, MSBuild.exe) which then beacons to DDNS infrastructure on common RAT ports (AsyncRAT, NetWire, WSH RAT, Parallax, Revenge RAT, njRAT). This parent\u2192child\u2192DDNS chain is distinctive to TA2541 because they consistently re-use this crypter family across aviation campaigns.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1055.012",
        "T1036.005",
        "T1027.002",
        "T1547.001",
        "T1053.005",
        "T1568",
        "T1573.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process) as parent_cmd values(Processes.process_path) as path from datamodel=Endpoint.Processes where Processes.process_name IN (\"RegSvcs.exe\",\"RegAsm.exe\",\"InstallUtil.exe\",\"AspNetCompiler.exe\",\"aspnet_compiler.exe\",\"MSBuild.exe\",\"jsc.exe\",\"vbc.exe\",\"csc.exe\") AND Processes.parent_process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"explorer.exe\") by host Processes.user Processes.process_name Processes.parent_process_name Processes.process Processes.parent_process \n| `drop_dm_object_name(Processes)` \n| join type=inner host [| tstats `summariesonly` count values(All_Traffic.dest) as ddns_dest values(All_Traffic.dest_port) as ports from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest IN (\"*duckdns.org\",\"*hopto.org\",\"*ddns.net\",\"*no-ip.com\",\"*sytes.net\",\"*zapto.org\",\"*linkpc.net\") AND All_Traffic.dest_port IN (1604,1177,2404,4404,5552,6606,7707,8808,2000,5000) by All_Traffic.src | rename All_Traffic.src as host] \n| join type=outer host [| tstats `summariesonly` count from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\CurrentVersion\\\\Run\\\\*\" by host Registry.registry_value_name Registry.registry_value_data | rename Registry.* as *] \n| eval mitre=\"T1055.012+T1036.005+T1568\", risk_score=90 \n| table firstTime lastTime host user parent_process_name process_name process ddns_dest ports registry_value_data risk_score",
      "kql": "let dotnetHosts = dynamic([\"RegSvcs.exe\",\"RegAsm.exe\",\"InstallUtil.exe\",\"AspNetCompiler.exe\",\"aspnet_compiler.exe\",\"MSBuild.exe\",\"jsc.exe\",\"vbc.exe\",\"csc.exe\"]);\nlet scriptParents = dynamic([\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"explorer.exe\"]);\nlet ddnsHosts = dynamic([\".duckdns.org\",\".hopto.org\",\".ddns.net\",\".no-ip.com\",\".sytes.net\",\".zapto.org\",\".linkpc.net\",\".myftp.org\",\".ydns.eu\"]);\nlet ratPorts = dynamic([1604,1177,2404,4404,5552,6606,7707,8808,2000,5000,8080,4444]);\nlet hollowed =\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName in~ (dotnetHosts)\n| where InitiatingProcessFileName in~ (scriptParents)\n| where FolderPath !startswith \"C:\\\\Windows\\\\Microsoft.NET\\\\\" or InitiatingProcessFileName !in~ (\"devenv.exe\",\"msdeploy.exe\")\n| project SuspectTime=Timestamp, DeviceId, DeviceName, AccountName, Parent=InitiatingProcessFileName, ParentCmd=InitiatingProcessCommandLine, Hollowed=FileName, HollowedCmd=ProcessCommandLine, HollowedPid=ProcessId;\nlet beacon =\nDeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where RemoteUrl has_any (ddnsHosts) or RemoteIPType == \"Public\"\n| where RemotePort in (ratPorts) or RemoteUrl has_any (ddnsHosts)\n| project NetTime=Timestamp, DeviceId, InitProc=InitiatingProcessFileName, InitPid=InitiatingProcessId, RemoteUrl, RemoteIP, RemotePort;\nlet persistence =\nDeviceRegistryEvents\n| where Timestamp > ago(7d)\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\")\n| project RegTime=Timestamp, DeviceId, RegistryValueName, RegistryValueData, RegInit=InitiatingProcessFileName;\nhollowed\n| join kind=inner beacon on DeviceId\n| where NetTime between (SuspectTime .. SuspectTime + 15m)\n| where InitProc =~ Hollowed or InitPid == HollowedPid\n| join kind=leftouter persistence on DeviceId\n| project SuspectTime, DeviceName, AccountName, Parent, ParentCmd, Hollowed, HollowedCmd, RemoteUrl, RemoteIP, RemotePort, RegistryValueName, RegistryValueData\n| order by SuspectTime desc",
      "rationale": "Across Cisco Talos (Operation Layover) and Proofpoint TA2541 reporting, the actor's crypter unpacks AsyncRAT/NetWire/Parallax/Revenge/WSH RAT into a hollowed signed .NET binary (RegSvcs, RegAsm, AspNetCompiler, MSBuild) launched from a script-host parent \u2014 a pattern almost never produced by legitimate developers in a non-dev workstation. Correlating that hollowing with DDNS C2 on canonical AsyncRA"
    }
  ]
}