# intel/.llm_actor_uc_cache/e8/e8dcf141f31abaa818a1d3a11ad5fc21d45ea376.json
{
  "name": "Sandworm",
  "ucs": [
    {
      "title": "[LLM] Sandworm wiper pre-detonation cascade: shadow-copy deletion + boot recovery tampering + event-log wipe within a 30 min window on a single ho",
      "description": "Hunts for the destructive-preparation sequence Sandworm consistently runs immediately before triggering a wiper or pseudo-ransomware payload (NotPetya, Olympic Destroyer, HermeticWiper, CaddyWiper, IsaacWiper, RansomBoggs, Prestige, Industroyer2 host components). The hallmark is multiple, distinct anti-recovery actions chained on one host inside a narrow window by a non-machine account.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1490",
        "T1485",
        "T1070.001",
        "T1489",
        "T1486"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count from datamodel=Endpoint.Processes where (   (Processes.process_name=vssadmin.exe AND (Processes.process=\"*delete*shadow*\" OR Processes.process=\"*shadow*delete*\"))   OR (Processes.process_name=wmic.exe AND Processes.process=\"*shadowcopy*delete*\")   OR (Processes.process_name=wbadmin.exe AND Processes.process=\"*delete*catalog*\")   OR (Processes.process_name=bcdedit.exe AND (Processes.process=\"*recoveryenabled*no*\" OR Processes.process=\"*bootstatuspolicy*ignoreallfailures*\"))   OR (Processes.process_name=wevtutil.exe AND Processes.process=\"*cl *\")   OR (Processes.process_name IN (powershell.exe,pwsh.exe) AND (Processes.process=\"*Clear-EventLog*\" OR Processes.process=\"*Remove-EventLog*\" OR Processes.process=\"*Get-WmiObject*Win32_Shadowcopy*\" OR Processes.process=\"*Remove-WmiObject*Win32_Shadowcopy*\"))   OR (Processes.process_name=fsutil.exe AND Processes.process=\"*usn*deletejournal*\")   OR (Processes.process_name=net.exe AND Processes.process=\"*stop*\" AND Processes.process IN (\"*MSSQL*\",\"*Veeam*\",\"*Backup*\",\"*Exchange*\",\"*MSExchange*\")) ) by Processes.dest Processes.user Processes.process_name Processes.process _time | `drop_dm_object_name(Processes)` | search NOT user IN (\"*$\",\"SYSTEM\",\"LOCAL SERVICE\",\"NETWORK SERVICE\") | bin _time span=30m | stats min(_time) as first_seen max(_time) as last_seen dc(process_name) as distinct_destructive_actions values(process_name) as actions values(process) as command_lines by dest user _time | where distinct_destructive_actions >= 2 | eval duration_minutes=round((last_seen-first_seen)/60,1) | sort - distinct_destructive_actions, dest",
      "kql": "let LookbackHours = 24h;\nlet WindowMinutes = 30;\nlet _destructive = DeviceProcessEvents\n    | where Timestamp > ago(LookbackHours)\n    | where AccountName !endswith \"$\"\n    | where AccountName !in~ (\"system\",\"local service\",\"network service\")\n    | extend ActionClass = case(\n        FileName =~ \"vssadmin.exe\" and ProcessCommandLine has \"delete\" and ProcessCommandLine has \"shadow\", \"ShadowDelete\",\n        FileName =~ \"wmic.exe\" and ProcessCommandLine has \"shadowcopy\" and ProcessCommandLine has \"delete\", \"ShadowDelete\",\n        FileName =~ \"wbadmin.exe\" and ProcessCommandLine has \"delete\" and ProcessCommandLine has \"catalog\", \"BackupCatalogDelete\",\n        FileName =~ \"bcdedit.exe\" and (ProcessCommandLine has_all (\"recoveryenabled\",\"no\") or ProcessCommandLine has \"bootstatuspolicy ignoreallfailures\"), \"BootRecoveryDisable\",\n        FileName =~ \"wevtutil.exe\" and ProcessCommandLine has \"cl \", \"EventLogClear\",\n        FileName in~ (\"powershell.exe\",\"pwsh.exe\") and ProcessCommandLine has_any (\"Clear-EventLog\",\"Remove-EventLog\",\"wevtutil cl\",\"Remove-WmiObject Win32_Shadowcopy\",\"Get-WmiObject Win32_Shadowcopy\"), \"PsRecoveryTamper\",\n        FileName =~ \"fsutil.exe\" and ProcessCommandLine has \"usn\" and ProcessCommandLine has \"deletejournal\", \"UsnJournalWipe\",\n        FileName =~ \"net.exe\" and ProcessCommandLine has \"stop\" and ProcessCommandLine has_any (\"MSSQL\",\"SQLAgent\",\"Veeam\",\"Backup\",\"MSExchange\",\"Exchange\"), \"BackupServiceStop\",\n        \"\")\n    | where isnotempty(ActionClass)\n    | project Timestamp, DeviceId, DeviceName, AccountName, ActionClass, FileName, ProcessCommandLine, InitiatingProcessFileName, InitiatingProcessCommandLine;\n_destructive\n| summarize ActionsSeen = make_set(ActionClass),\n            DistinctActions = dcount(ActionClass),\n            CommandLines = make_set(ProcessCommandLine, 25),\n            Parents = make_set(InitiatingProcessFileName, 10),\n            FirstSeen = min(Timestamp),\n            LastSeen = max(Timestamp)\n            by DeviceName, DeviceId, AccountName, bin(Timestamp, WindowMinutes * 1m)\n| where DistinctActions >= 2\n| extend WindowMinutes = datetime_diff('minute', LastSeen, FirstSeen)\n| project FirstSeen, LastSeen, WindowMinutes, DeviceName, AccountName, DistinctActions, ActionsSeen, Parents, CommandLines\n| order by DistinctActions desc, FirstSeen desc",
      "sentinel_kql": "let LookbackHours = 24h;\nlet WindowMinutes = 30;\nlet _destructive = SecurityEvent\n    | where TimeGenerated > ago(LookbackHours)\n    | where EventID == 4688\n    | where AccountType == \"User\"\n    | where Account !endswith \"$\"\n    | extend Img = tolower(NewProcessName), Cmd = tolower(CommandLine)\n    | extend ActionClass = case(\n        Img endswith \"\\\\vssadmin.exe\" and Cmd has \"delete\" and Cmd has \"shadow\", \"ShadowDelete\",\n        Img endswith \"\\\\wmic.exe\" and Cmd has \"shadowcopy\" and Cmd has \"delete\", \"ShadowDelete\",\n        Img endswith \"\\\\wbadmin.exe\" and Cmd has \"delete\" and Cmd has \"catalog\", \"BackupCatalogDelete\",\n        Img endswith \"\\\\bcdedit.exe\" and (Cmd has_all (\"recoveryenabled\",\"no\") or Cmd has \"bootstatuspolicy ignoreallfailures\"), \"BootRecoveryDisable\",\n        Img endswith \"\\\\wevtutil.exe\" and Cmd has \"cl \", \"EventLogClear\",\n        (Img endswith \"\\\\powershell.exe\" or Img endswith \"\\\\pwsh.exe\") and Cmd has_any (\"clear-eventlog\",\"remove-eventlog\",\"wevtutil cl\",\"remove-wmiobject win32_shadowcopy\"), \"PsRecoveryTamper\",\n        Img endswith \"\\\\fsutil.exe\" and Cmd has \"usn\" and Cmd has \"deletejournal\", \"UsnJournalWipe\",\n        Img endswith \"\\\\net.exe\" and Cmd has \"stop\" and Cmd has_any (\"mssql\",\"sqlagent\",\"veeam\",\"backup\",\"msexchange\",\"exchange\"), \"BackupServiceStop\",\n        \"\")\n    | where isnotempty(ActionClass);\n_destructive\n| summarize ActionsSeen = make_set(ActionClass),\n            DistinctActions = dcount(ActionClass),\n            CommandLines = make_set(CommandLine, 25),\n            Parents = make_set(ParentProcessName, 10),\n            FirstSeen = min(TimeGenerated),\n            LastSeen = max(TimeGenerated)\n            by Computer, Account, bin(TimeGenerated, WindowMinutes * 1m)\n| where DistinctActions >= 2\n| extend WindowMinutes = datetime_diff('minute', LastSeen, FirstSeen)\n| project FirstSeen, LastSeen, WindowMinutes, Computer, Account, DistinctActions, ActionsSeen, Parents, CommandLines\n| order by DistinctActions desc, FirstSeen desc",
      "sigma_yaml": "title: Sandworm Wiper Pre-Detonation Anti-Recovery Command\nid: a8b1e6e0-3c4c-4d2f-9bb7-sandworm-precursor\nstatus: experimental\ndescription: Detects individual destructive-preparation actions characteristic of Sandworm wiper deployments (NotPetya, HermeticWiper, CaddyWiper, IsaacWiper, RansomBoggs, Prestige, Industroyer2). Correlate multiple hits on the same host within 30 min for high-fidelity alerting.\nauthor: Detection Engineering\nreferences:\n  - https://attack.mitre.org/groups/G0034/\n  - https://www.cisa.gov/news-events/cybersecurity-advisories/aa22-057a\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  shadow_delete:\n    Image|endswith:\n      - '\\vssadmin.exe'\n      - '\\wmic.exe'\n    CommandLine|contains|all:\n      - 'shadow'\n      - 'delete'\n  catalog_delete:\n    Image|endswith: '\\wbadmin.exe'\n    CommandLine|contains|all:\n      - 'delete'\n      - 'catalog'\n  boot_recovery_off:\n    Image|endswith: '\\bcdedit.exe'\n    CommandLine|contains:\n      - 'recoveryenabled no'\n      - 'bootstatuspolicy ignoreallfailures'\n  event_log_wipe:\n    Image|endswith: '\\wevtutil.exe'\n    CommandLine|contains: 'cl '\n  ps_recovery_tamper:\n    Image|endswith:\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n    CommandLine|contains:\n      - 'Clear-EventLog'\n      - 'Remove-EventLog'\n      - 'Remove-WmiObject Win32_Shadowcopy'\n  usn_wipe:\n    Image|endswith: '\\fsutil.exe'\n    CommandLine|contains|all:\n      - 'usn'\n      - 'deletejournal'\n  filter_machine_acct:\n    User|endswith: '$'\n  condition: 1 of (shadow_delete, catalog_delete, boot_recovery_off, event_log_wipe, ps_recovery_tamper, usn_wipe) and not filter_machine_acct\nfalsepositives:\n  - Backup software life-cycle scripts (Veeam, BackupExec) running as service accounts\n  - Authorised disaster-recovery testing\nlevel: high\ntags:\n  - attack.impact\n  - attack.t1490\n  - attack.t1485\n  - attack.t1070.001\n  - attack.g0034\n",
      "rationale": "Sandworm's destructive operations from BlackEnergy through Industroyer2 share a deliberate anti-forensic / anti-recovery cascade run seconds before the wiper or NotPetya-style pseudo-ransomware fires: vssadmin/wmic shadow-copy nuke, wbadmin catalog delete, bcdedit boot-recovery off, wevtutil/Clear-EventLog wipe, and a service-stop sweep to free backup/Exchange/SQL files for overwrite. Other actors",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] Sandworm centralised wiper push: identical non-system binary executed near-simultaneously across N hosts via SchTasks/WMI/PsExec/GPO from a ",
      "description": "Hunts the deployment fingerprint Sandworm reuses across NotPetya (M.E.Doc + Mimikatz + PsExec/WMIC), Industroyer2 (scheduled tasks pushed from a compromised DC), Prestige (Group Policy + remote scheduled task), and RansomBoggs (PowerShell DC distribution): the same SHA256 starts executing on many hosts within minutes, parented by a remote-execution mechanism, from a non-Microsoft path. Legitimate enterprise software deployment does not deliver to dozens of hosts within a sub-10-minute window from \\Windows\\ or \\PerfLogs\\.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1053.005",
        "T1047",
        "T1021.002",
        "T1072",
        "T1570",
        "T1486",
        "T1485"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count min(_time) as first_seen max(_time) as last_seen values(Processes.dest) as hosts values(Processes.user) as users values(Processes.parent_process_name) as parents values(Processes.process) as sample_cmd from datamodel=Endpoint.Processes where (Processes.parent_process_name IN (schtasks.exe,wmic.exe,services.exe,svchost.exe,psexec.exe,psexec64.exe,paexec.exe,wmiprvse.exe,taskeng.exe) OR Processes.grandparent_process_name IN (schtasks.exe,wmic.exe,services.exe,wmiprvse.exe)) AND (Processes.process_path=\"*\\\\Windows\\\\*\" OR Processes.process_path=\"*\\\\PerfLogs\\\\*\" OR Processes.process_path=\"*\\\\ProgramData\\\\*\" OR Processes.process_path=\"*\\\\Users\\\\Public\\\\*\" OR Processes.process_path=\"*\\\\Temp\\\\*\") AND NOT (Processes.process_path=\"*\\\\Windows\\\\System32\\\\*\" OR Processes.process_path=\"*\\\\Windows\\\\SysWOW64\\\\*\" OR Processes.process_path=\"*\\\\Windows\\\\WinSxS\\\\*\" OR Processes.process_path=\"*\\\\Windows\\\\servicing\\\\*\" OR Processes.process_path=\"*\\\\Windows\\\\SoftwareDistribution\\\\*\") AND NOT Processes.process_name IN (svchost.exe,conhost.exe,wermgr.exe,WerFault.exe,MpCmdRun.exe,MsMpEng.exe,SenseIR.exe,MoUsoCoreWorker.exe) by Processes.process_hash Processes.process_name | `drop_dm_object_name(Processes)` | eval host_count=mvcount(hosts), duration_seconds=last_seen-first_seen | where host_count>=5 AND duration_seconds<=600 | eval rate_hosts_per_min=round(host_count/((duration_seconds/60)+0.1),1) | sort - host_count, duration_seconds",
      "kql": "let LookbackHours = 6h;\nlet WindowMinutes = 10;\nlet MinHosts = 5;\nlet _remote_launch_parents = dynamic([\"schtasks.exe\",\"wmic.exe\",\"wmiprvse.exe\",\"services.exe\",\"psexec.exe\",\"psexec64.exe\",\"paexec.exe\",\"taskeng.exe\",\"svchost.exe\"]);\nlet _benign_paths = dynamic([@\"\\Windows\\System32\\\", @\"\\Windows\\SysWOW64\\\", @\"\\Windows\\WinSxS\\\", @\"\\Windows\\servicing\\\", @\"\\Windows\\SoftwareDistribution\\\", @\"\\Windows\\Microsoft.NET\\\", @\"\\Program Files\\Windows Defender\\\", @\"\\Program Files (x86)\\Microsoft\\\"]);\nlet _benign_names = dynamic([\"svchost.exe\",\"conhost.exe\",\"wermgr.exe\",\"werfault.exe\",\"mpcmdrun.exe\",\"msmpeng.exe\",\"senseir.exe\",\"mousocoreworker.exe\",\"trustedinstaller.exe\",\"tiworker.exe\",\"smartscreen.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(LookbackHours)\n| where InitiatingProcessFileName in~ (_remote_launch_parents)\n     or InitiatingProcessParentFileName in~ (_remote_launch_parents)\n| where FileName !in~ (_benign_names)\n| where FolderPath has_any (@\"\\Windows\\\", @\"\\PerfLogs\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\Temp\\\", @\"\\AppData\\Local\\Temp\\\")\n| where not(FolderPath has_any (_benign_paths))\n| where isnotempty(SHA256)\n| summarize HostCount = dcount(DeviceId),\n            Hosts = make_set(DeviceName, 100),\n            Users = make_set(AccountName, 20),\n            Parents = make_set(InitiatingProcessFileName, 10),\n            SampleCmd = any(ProcessCommandLine),\n            SamplePath = any(FolderPath),\n            FirstSeen = min(Timestamp),\n            LastSeen = max(Timestamp)\n            by SHA256, FileName\n| where HostCount >= MinHosts\n| extend SpreadMinutes = datetime_diff('minute', LastSeen, FirstSeen)\n| where SpreadMinutes <= WindowMinutes\n| extend HostsPerMinute = round(toreal(HostCount) / (toreal(SpreadMinutes) + 0.5), 1)\n| project FirstSeen, LastSeen, SpreadMinutes, HostCount, HostsPerMinute, FileName, SHA256, SamplePath, Parents, Users, SampleCmd, Hosts\n| order by HostCount desc, SpreadMinutes asc",
      "sentinel_kql": "let LookbackHours = 6h;\nlet WindowMinutes = 10;\nlet MinHosts = 5;\nlet _remote_launch_parents = dynamic([\"schtasks.exe\",\"wmic.exe\",\"wmiprvse.exe\",\"services.exe\",\"psexec.exe\",\"psexec64.exe\",\"paexec.exe\",\"taskeng.exe\",\"svchost.exe\"]);\nlet _benign_names = dynamic([\"svchost.exe\",\"conhost.exe\",\"wermgr.exe\",\"werfault.exe\",\"mpcmdrun.exe\",\"msmpeng.exe\",\"trustedinstaller.exe\",\"tiworker.exe\",\"smartscreen.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(LookbackHours)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| extend ParentLeaf = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1])),\n         ChildLeaf  = tolower(tostring(split(NewProcessName, \"\\\\\")[-1])),\n         ImgLow     = tolower(NewProcessName)\n| where ParentLeaf in (_remote_launch_parents)\n| where ChildLeaf !in (_benign_names)\n| where ImgLow has_any (\"\\\\windows\\\\\", \"\\\\perflogs\\\\\", \"\\\\programdata\\\\\", \"\\\\users\\\\public\\\\\", \"\\\\temp\\\\\", \"\\\\appdata\\\\local\\\\temp\\\\\")\n| where not(ImgLow has_any (\"\\\\windows\\\\system32\\\\\", \"\\\\windows\\\\syswow64\\\\\", \"\\\\windows\\\\winsxs\\\\\", \"\\\\windows\\\\servicing\\\\\", \"\\\\windows\\\\softwaredistribution\\\\\", \"\\\\windows\\\\microsoft.net\\\\\"))\n| where isnotempty(FileHash)\n| summarize HostCount = dcount(Computer),\n            Hosts = make_set(Computer, 100),\n            Users = make_set(Account, 20),\n            Parents = make_set(ParentLeaf, 10),\n            SampleCmd = any(CommandLine),\n            SamplePath = any(NewProcessName),\n            FirstSeen = min(TimeGenerated),\n            LastSeen = max(TimeGenerated)\n            by FileHash, ChildLeaf\n| where HostCount >= MinHosts\n| extend SpreadMinutes = datetime_diff('minute', LastSeen, FirstSeen)\n| where SpreadMinutes <= WindowMinutes\n| extend HostsPerMinute = round(toreal(HostCount) / (toreal(SpreadMinutes) + 0.5), 1)\n| project FirstSeen, LastSeen, SpreadMinutes, HostCount, HostsPerMinute, ChildLeaf, FileHash, SamplePath, Parents, Users, SampleCmd, Hosts\n| order by HostCount desc, SpreadMinutes asc",
      "sigma_yaml": "",
      "rationale": "Sandworm's distinguishing operational pattern is centralised mass deployment of the destructive payload \u2014 not the wiper code itself. NotPetya rode M.E.Doc updates + WMIC/PsExec laterally, Prestige (Microsoft IRIS-T2) was pushed via Active Directory Group Policy and remote scheduled tasks, Industroyer2 was distributed by scheduled tasks dropped from a compromised DC, and RansomBoggs was launched vi"
    }
  ]
}