# intel/.llm_actor_uc_cache/68/682ece814f199590cbff65ddc30437d94a5aed1d.json
{
  "name": "Winter Vivern",
  "ucs": [
    {
      "title": "[LLM] Winter Vivern (TA473 / UAC-0114) macro-borne PowerShell creating scheduled task with embedded URL/encoded payload",
      "description": "Hunts the canonical Winter Vivern infection chain reported by SentinelLabs, CERT-UA and Proofpoint: a weaponised Office document or HTA spawns a shell which immediately calls schtasks.exe /create, registering a task that re-launches PowerShell with an attacker URL or base64 blob. This Office\u2192shell\u2192schtasks sequence inside one process tree is a near-fingerprint of the group's APERETIF / Aperetif PowerShell loader and the Polish/Ukrainian phishing waves CERT-UA tracks as UAC-0114.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.001",
        "T1059.001",
        "T1059.003",
        "T1053.005",
        "T1105",
        "T1071.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.process_name=schtasks.exe Processes.process=\"*/create*\" by _time host Processes.user Processes.parent_process_id Processes.process | rename Processes.* as * | where match(process, \"(?i)(powershell|pwsh|invoke-webrequest|iex|downloadstring|curl|http://|https://|-enc|-encodedcommand|mshta|certutil)\") | rename process as schtasks_cmd parent_process_id as shell_pid | join type=inner host shell_pid [ | tstats `summariesonly` count from datamodel=Endpoint.Processes where Processes.process_name IN (powershell.exe, pwsh.exe, cmd.exe) Processes.parent_process_name IN (winword.exe, excel.exe, powerpnt.exe, outlook.exe, wscript.exe, cscript.exe, mshta.exe, hh.exe) by host Processes.process_id Processes.parent_process_name Processes.process | rename Processes.* as *, process_id as shell_pid, parent_process_name as office_parent, process as shell_cmd ] | where (lastTime - firstTime) < 600 | table _time host user office_parent shell_cmd schtasks_cmd | sort - _time",
      "kql": "let CorrWin = 10m;\nlet OfficeAndScriptHosts = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"hh.exe\"]);\n// Stage A \u2014 Office / script-host spawning a Windows shell (macro detonation)\nlet StageA = DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where InitiatingProcessFileName in~ (OfficeAndScriptHosts)\n    | where FileName in~ (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\")\n    | where AccountName !endswith \"$\"\n    | project StageATime = Timestamp, DeviceId, DeviceName, AccountName,\n              OfficeParent = InitiatingProcessFileName,\n              ShellCmd = ProcessCommandLine,\n              ShellPid = ProcessId;\n// Stage B \u2014 that shell creates a scheduled task whose action carries a URL / encoded payload\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has \"/create\"\n| where AccountName !endswith \"$\"\n| join kind=inner StageA on DeviceId, $left.InitiatingProcessId == $right.ShellPid\n| where Timestamp between (StageATime .. StageATime + CorrWin)\n// TA473 / UAC-0114 schtasks lines re-launch PS with a URL or base64 \u2014 the give-away.\n| where ProcessCommandLine has_any (\n    \"powershell\",\"pwsh\",\"Invoke-WebRequest\",\"IEX\",\"DownloadString\",\n    \"curl\",\"certutil\",\"http://\",\"https://\",\"-enc\",\"-EncodedCommand\",\"mshta\"\n  )\n| project StageATime, SchTime = Timestamp, DeviceName, AccountName,\n          OfficeParent, ShellCmd, SchTaskCmd = ProcessCommandLine\n| order by SchTime desc",
      "sentinel_kql": "let CorrWin = 10m;\nlet OfficeAndScriptHosts = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"hh.exe\"]);\n// Stage A \u2014 Office / script-host spawns a shell (4688 child of macro host)\nlet StageA = SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ParentLeaf = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1])),\n             ChildLeaf  = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n    | where ParentLeaf in (OfficeAndScriptHosts)\n    | where ChildLeaf in (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\")\n    | project StageATime = TimeGenerated, Computer, Account,\n              OfficeParent = ParentLeaf,\n              ShellCmd = CommandLine,\n              ShellPid = NewProcessId;\n// Stage B \u2014 that shell calls schtasks /create with a URL or encoded blob\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| extend NewLeaf = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n| where NewLeaf == \"schtasks.exe\"\n| where CommandLine has \"/create\"\n| join kind=inner StageA on Computer, $left.ProcessId == $right.ShellPid\n| where TimeGenerated between (StageATime .. StageATime + CorrWin)\n| where CommandLine has_any (\n    \"powershell\",\"pwsh\",\"Invoke-WebRequest\",\"IEX\",\"DownloadString\",\n    \"curl\",\"certutil\",\"http://\",\"https://\",\"-enc\",\"-EncodedCommand\",\"mshta\"\n  )\n| project StageATime, SchTime = TimeGenerated, Computer, Account,\n          OfficeParent, ShellCmd, SchTaskCmd = CommandLine\n| order by SchTime desc",
      "sigma_yaml": "",
      "rationale": "This is not a generic 'Office spawns PowerShell' rule \u2014 it requires (1) a script-host macro/HTA detonation as parent, (2) a child shell, AND (3) that exact shell child instance launching schtasks /create within 10 minutes carrying a URL or -EncodedCommand. SentinelOne's APERETIF write-up, CERT-UA's UAC-0114 advisories and Proofpoint's TA473 reporting all describe this same three-step payload-stage"
    },
    {
      "title": "[LLM] Winter Vivern lookalike-domain egress from PowerShell / LOLBins (gov / MFA / OSCE typosquat C2)",
      "description": "Surface PowerShell, curl, certutil, mshta or BITS reaching out to domains that impersonate European government, ministry-of-foreign-affairs or NGO branding (e.g. ocse[.]com, poczta-mil[.]com, mod-gov[.]uk, mfa-gov[.]net) on non-government TLDs. Winter Vivern's 2021-2024 campaigns are catalogued (DomainTools, SentinelLabs, ESET) for registering exactly this style of typosquat as both staging and C2 infrastructure for their APERETIF/AMETHYST loaders.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "medium",
      "techs": [
        "T1583.001",
        "T1071.001",
        "T1105",
        "T1036.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime values(All_Traffic.dest) as dest_ip values(Web.url) as urls from datamodel=Web.Web where Web.user_agent_length>0 (Web.process IN (\"*\\\\powershell.exe\",\"*\\\\pwsh.exe\",\"*\\\\curl.exe\",\"*\\\\certutil.exe\",\"*\\\\mshta.exe\",\"*\\\\bitsadmin.exe\",\"*\\\\wscript.exe\",\"*\\\\cscript.exe\")) by host Web.user Web.process Web.url Web.dest | rename Web.* as * | rex field=url \"https?://(?<fqdn>[^/]+)/\" | where match(fqdn, \"(?i)(gov|mfa|mod|mil|ministry|embassy|council|europa|ocse|osce|poczta|ukraine|polska|poland)\") AND NOT match(fqdn, \"(?i)(\\.gov($|\\.[a-z]{2,3}$)|\\.mil$|\\.mod\\.uk$|europa\\.eu$)\") | stats min(_time) as firstSeen max(_time) as lastSeen values(host) as hosts dc(host) as host_count values(user) as users values(process) as tools values(url) as sample_urls by fqdn | sort - lastSeen",
      "kql": "let GovKeywords = dynamic([\"gov\",\"mfa\",\"mod\",\"mil\",\"ministry\",\"embassy\",\"council\",\"europa\",\"ocse\",\"osce\",\"poczta\",\"ukraine\",\"polska\",\"poland\"]);\nlet LegitGovTLDPattern = @\"(?i)(\\.gov($|\\.[a-z]{2,3}($|/))|\\.mil($|/)|\\.mod\\.uk($|/)|europa\\.eu($|/))\";\nlet LolBins = dynamic([\"powershell.exe\",\"pwsh.exe\",\"curl.exe\",\"certutil.exe\",\"mshta.exe\",\"bitsadmin.exe\",\"wscript.exe\",\"cscript.exe\",\"rundll32.exe\",\"regsvr32.exe\"]);\nDeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (LolBins)\n| where RemoteIPType == \"Public\"\n| where isnotempty(RemoteUrl)\n| extend Fqdn = tolower(tostring(extract(@\"https?://([^/]+)/?\", 1, RemoteUrl)))\n| where isnotempty(Fqdn)\n| where Fqdn has_any (GovKeywords)\n| where not(RemoteUrl matches regex LegitGovTLDPattern)\n// Anti-baseline: brand-new contact in the last 24h\n| summarize FirstSeen = min(Timestamp), LastSeen = max(Timestamp),\n            Hosts = make_set(DeviceName, 25),\n            HostCount = dcount(DeviceName),\n            Tools = make_set(InitiatingProcessFileName, 10),\n            SampleUrl = any(RemoteUrl),\n            SampleCmd = any(InitiatingProcessCommandLine)\n            by Fqdn, RemoteIP\n| where FirstSeen > ago(24h)\n| order by FirstSeen desc",
      "sentinel_kql": "let GovKeywords = dynamic([\"gov\",\"mfa\",\"mod\",\"mil\",\"ministry\",\"embassy\",\"council\",\"europa\",\"ocse\",\"osce\",\"poczta\",\"ukraine\",\"polska\",\"poland\"]);\nlet LegitGovTLDPattern = @\"(?i)(\\.gov($|\\.[a-z]{2,3}($|/))|\\.mil($|/)|\\.mod\\.uk($|/)|europa\\.eu($|/))\";\nlet LolBins = dynamic([\"powershell.exe\",\"pwsh.exe\",\"curl.exe\",\"certutil.exe\",\"mshta.exe\",\"bitsadmin.exe\",\"wscript.exe\",\"cscript.exe\",\"rundll32.exe\",\"regsvr32.exe\"]);\nImWebSession\n| where TimeGenerated > ago(7d)\n| extend ToolLeaf = tolower(tostring(split(ActingProcessName, \"\\\\\")[-1]))\n| where ToolLeaf in (LolBins)\n| where isnotempty(Url)\n| extend Fqdn = tolower(coalesce(HttpHostName, tostring(extract(@\"https?://([^/]+)/?\", 1, Url))))\n| where isnotempty(Fqdn)\n| where Fqdn has_any (GovKeywords)\n| where not(Url matches regex LegitGovTLDPattern)\n| summarize FirstSeen = min(TimeGenerated), LastSeen = max(TimeGenerated),\n            Hosts = make_set(DvcHostname, 25),\n            HostCount = dcount(DvcHostname),\n            Tools = make_set(ToolLeaf, 10),\n            SampleUrl = any(Url),\n            SampleCmd = any(ActingProcessCommandLine)\n            by Fqdn, DstIpAddr\n| where FirstSeen > ago(24h)\n| order by FirstSeen desc",
      "sigma_yaml": "",
      "rationale": "Winter Vivern's distinguishing operational tradecraft (per DomainTools 2021, SentinelLabs 2023, ESET 2023 reporting) is procuring lookalike domains that imitate the exact victim ministries, embassies, OSCE, the European Council, and Polish/Ukrainian government mail systems \u2014 then using them as both phishing pretexts and as PowerShell-loader C2. Constraining the source process to Windows LOLBins an",
      "_sentinel_field_issues": [
        {
          "kind": "wrong_table",
          "field": "ActingProcessName",
          "tables_in_scope": [
            "ImWebSession"
          ],
          "exists_on": [
            "ImProcessCreate",
            "ImFileEvent",
            "ImRegistryEvent"
          ],
          "message": "`ActingProcessName` is not a column of ImWebSession; it does exist on ImProcessCreate, ImFileEvent, ImRegistryEvent."
        },
        {
          "kind": "wrong_table",
          "field": "ActingProcessCommandLine",
          "tables_in_scope": [
            "ImWebSession"
          ],
          "exists_on": [
            "ImProcessCreate",
            "ImFileEvent",
            "ImRegistryEvent"
          ],
          "message": "`ActingProcessCommandLine` is not a column of ImWebSession; it does exist on ImProcessCreate, ImFileEvent, ImRegistryEvent."
        }
      ]
    }
  ]
}