# intel/.llm_actor_uc_cache/e3/e3c64d5bb80e2220052561617dbe293d64b412ea.json
{
  "name": "Silk Typhoon",
  "ucs": [
    {
      "title": "[LLM] Silk Typhoon (HAFNIUM) Exchange w3wp.exe spawning interpreters/recon \u2014 ProxyLogon-style web shell post-exploitation",
      "description": "Hunts the canonical Silk Typhoon / HAFNIUM tradecraft of exploiting on-prem Exchange (ProxyLogon CVE-2021-26855 chain and successors) and using the IIS Exchange app-pool worker process w3wp.exe to spawn cmd/PowerShell, recon LOLBins, or procdump as a child of the dropped China-Chopper-style .aspx web shell. The parent w3wp.exe with an Exchange app-pool command line spawning these binaries is essentially never benign.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "exploit",
      "conf": "high",
      "techs": [
        "T1190",
        "T1505.003",
        "T1059.001",
        "T1059.003",
        "T1003.001",
        "T1082"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.parent_process_name=w3wp.exe Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\",\"cscript.exe\",\"wscript.exe\",\"net.exe\",\"net1.exe\",\"whoami.exe\",\"systeminfo.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"quser.exe\",\"procdump.exe\",\"procdump64.exe\",\"rundll32.exe\",\"certutil.exe\",\"bitsadmin.exe\") by Processes.dest Processes.user Processes.parent_process Processes.parent_process_id Processes.process Processes.process_name Processes.process_id Processes.process_hash | `drop_dm_object_name(Processes)` | where like(parent_process,\"%MSExchange%AppPool%\") OR like(parent_process,\"%w3wp.exe -ap %Exchange%\") | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "// Silk Typhoon (HAFNIUM) \u2014 Exchange w3wp post-exploit child process\nlet ExchangeAppPools = dynamic([\"MSExchangeOWAAppPool\",\"MSExchangeECPAppPool\",\"MSExchangeAutodiscoverAppPool\",\"MSExchangeServicesAppPool\",\"MSExchangeRpcProxyAppPool\",\"MSExchangeRestAppPool\",\"MSExchangeMapiFrontEndAppPool\",\"MSExchangeMapiMailboxAppPool\",\"MSExchangePowerShellAppPool\",\"MSExchangePowerShellFrontEndAppPool\"]);\nlet HafniumChildBins = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\",\"cscript.exe\",\"wscript.exe\",\"net.exe\",\"net1.exe\",\"whoami.exe\",\"systeminfo.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"quser.exe\",\"nltest.exe\",\"procdump.exe\",\"procdump64.exe\",\"rundll32.exe\",\"certutil.exe\",\"bitsadmin.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(30d)\n| where InitiatingProcessFileName =~ \"w3wp.exe\"\n| where InitiatingProcessCommandLine has_any (ExchangeAppPools)\n| where FileName in~ (HafniumChildBins)\n| where AccountName !endswith \"$\"   // exclude machine-account scheduled work\n| project Timestamp, DeviceName, AccountName,\n          ParentCmd = InitiatingProcessCommandLine,\n          Child = FileName, ChildCmd = ProcessCommandLine,\n          ChildSHA256 = SHA256, FolderPath\n| order by Timestamp desc",
      "sentinel_kql": "// Silk Typhoon (HAFNIUM) \u2014 Exchange w3wp post-exploit child process (Sentinel SecurityEvent 4688)\nlet ExchangeAppPools = dynamic([\"MSExchangeOWAAppPool\",\"MSExchangeECPAppPool\",\"MSExchangeAutodiscoverAppPool\",\"MSExchangeServicesAppPool\",\"MSExchangeRpcProxyAppPool\",\"MSExchangeRestAppPool\",\"MSExchangeMapiFrontEndAppPool\",\"MSExchangeMapiMailboxAppPool\",\"MSExchangePowerShellAppPool\",\"MSExchangePowerShellFrontEndAppPool\"]);\nlet HafniumChildPaths = dynamic([\"\\\\cmd.exe\",\"\\\\powershell.exe\",\"\\\\pwsh.exe\",\"\\\\mshta.exe\",\"\\\\cscript.exe\",\"\\\\wscript.exe\",\"\\\\net.exe\",\"\\\\net1.exe\",\"\\\\whoami.exe\",\"\\\\systeminfo.exe\",\"\\\\ipconfig.exe\",\"\\\\tasklist.exe\",\"\\\\quser.exe\",\"\\\\nltest.exe\",\"\\\\procdump.exe\",\"\\\\procdump64.exe\",\"\\\\rundll32.exe\",\"\\\\certutil.exe\",\"\\\\bitsadmin.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(30d)\n| where EventID == 4688\n| where ParentProcessName endswith \"\\\\w3wp.exe\"\n| where NewProcessName has_any (HafniumChildPaths)\n| where Account !endswith \"$\"\n| extend AppPoolHint = extract(@\"-ap\\s+\\\"?([^\\\"\\s]+)\\\"?\", 1, CommandLine)\n| project TimeGenerated, Computer, Account, ParentProcessName,\n          NewProcessName, CommandLine, TokenElevationType, AppPoolHint\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Silk Typhoon (HAFNIUM) Exchange Web Shell Post-Exploitation Child Process\nid: 6f5d9c1e-3e6d-4b6a-9d5b-1a3a3c1f8d2c\nstatus: experimental\ndescription: |\n  Detects w3wp.exe (Exchange IIS app pool) spawning interpreters, recon LOLBins, or credential-dumping tooling \u2014 the canonical post-exploitation pattern of Silk Typhoon / HAFNIUM (G0125) following ProxyLogon-class Exchange exploitation and China-Chopper-style .aspx web shell drops.\nreferences:\n  - https://www.microsoft.com/en-us/security/blog/2021/03/02/hafnium-targeting-exchange-servers/\n  - https://www.microsoft.com/en-us/security/blog/2025/03/05/silk-typhoon-targeting-it-supply-chain/\n  - https://attack.mitre.org/groups/G0125/\nauthor: SOC Detection Engineering\ndate: 2026/05/10\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent:\n    ParentImage|endswith: '\\w3wp.exe'\n  selection_child:\n    Image|endswith:\n      - '\\cmd.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\mshta.exe'\n      - '\\cscript.exe'\n      - '\\wscript.exe'\n      - '\\net.exe'\n      - '\\net1.exe'\n      - '\\whoami.exe'\n      - '\\systeminfo.exe'\n      - '\\ipconfig.exe'\n      - '\\tasklist.exe'\n      - '\\quser.exe'\n      - '\\nltest.exe'\n      - '\\procdump.exe'\n      - '\\procdump64.exe'\n      - '\\rundll32.exe'\n      - '\\certutil.exe'\n      - '\\bitsadmin.exe'\n  filter_user:\n    User|endswith: '$'\n  condition: selection_parent and selection_child and not filter_user\nfields:\n  - ComputerName\n  - User\n  - Image\n  - CommandLine\n  - ParentImage\n  - ParentCommandLine\nfalsepositives:\n  - Legitimate Exchange administration scripts run inside an OWA/ECP request context (rare; should be allow-listed by hash + admin account)\nlevel: critical\ntags:\n  - attack.initial_access\n  - attack.t1190\n  - attack.persistence\n  - attack.t1505.003\n  - attack.execution\n  - attack.t1059.001\n  - attack.t1059.003\n  - attack.g0125",
      "rationale": "Silk Typhoon / HAFNIUM (G0125, Operation Exchange Marauder) is the named actor that mass-exploited on-prem Exchange via ProxyLogon (CVE-2021-26855/26857/26858/27065) in early 2021 and continues to chase Exchange/IT-edge bugs through 2024\u20132026. Their distinctive post-exploitation pattern, documented by MSTIC/DART, is to drop China-Chopper / Antsword .aspx shells under Exchange virtual directories a"
    },
    {
      "title": "[LLM] Silk Typhoon (HAFNIUM) Tarrask hidden scheduled-task persistence \u2014 TaskCache\\Tree SD value deletion",
      "description": "Detects the Tarrask malware persistence trick \u2014 attributed by Microsoft directly to HAFNIUM \u2014 where the actor creates a Windows scheduled task and then deletes the task's Security Descriptor (`SD`) registry value under `HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree\\<TaskName>\\SD`. With the SD value missing, the task is invisible to schtasks.exe, the Task Scheduler MMC, and `schtasks /query`, but still executes \u2014 a near-unique HAFNIUM defence-evasion signature.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1053.005",
        "T1562.001",
        "T1070.009",
        "T1027"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Registry where Registry.action=deleted Registry.registry_path=\"*\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Schedule\\\\TaskCache\\\\Tree\\\\*\" Registry.registry_value_name=SD by Registry.dest Registry.user Registry.process_guid Registry.process_id Registry.registry_path Registry.registry_value_name | `drop_dm_object_name(Registry)` | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "// Silk Typhoon (HAFNIUM) \u2014 Tarrask hidden scheduled task: SD value deleted from TaskCache\\Tree\nDeviceRegistryEvents\n| where Timestamp > ago(30d)\n| where ActionType in (\"RegistryValueDeleted\",\"RegistryKeyDeleted\")\n| where RegistryKey has @\"\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree\\\"\n| where RegistryValueName =~ \"SD\"\n   or RegistryKey endswith @\"\\SD\"\n| where InitiatingProcessFileName !in~ (\"trustedinstaller.exe\",\"msiexec.exe\")  // legit servicing rarely touches this\n| project Timestamp, DeviceName,\n          TaskTreePath = RegistryKey,\n          ValueDeleted = RegistryValueName,\n          Actor = InitiatingProcessAccountName,\n          ActorProcess = InitiatingProcessFileName,\n          ActorCmd = InitiatingProcessCommandLine,\n          ActorParent = InitiatingProcessParentFileName,\n          ActorIntegrity = InitiatingProcessIntegrityLevel\n| order by Timestamp desc",
      "sentinel_kql": "// Silk Typhoon (HAFNIUM) \u2014 Tarrask hidden scheduled task: SD value deleted from TaskCache\\Tree (ASIM)\nImRegistryEvent\n| where TimeGenerated > ago(30d)\n| where EventType in (\"RegistryValueDeleted\",\"RegistryKeyDeleted\",\"DeleteValue\")\n| where RegistryKey contains @\"\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree\\\"\n| where RegistryValue =~ \"SD\"\n   or RegistryKey endswith @\"\\SD\"\n| where ActingProcessName !endswith \"\\\\trustedinstaller.exe\"\n   and ActingProcessName !endswith \"\\\\msiexec.exe\"\n| project TimeGenerated, DvcHostname,\n          TaskTreePath = RegistryKey,\n          ValueDeleted = RegistryValue,\n          Actor = ActorUsername,\n          ActorProcess = ActingProcessName,\n          ActorCmd = ActingProcessCommandLine\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Silk Typhoon (HAFNIUM) Tarrask Hidden Scheduled Task SD-Value Deletion\nid: 8d4a5b21-7c1e-4f3b-9b22-a7c1e9d4c2f1\nstatus: experimental\ndescription: |\n  Detects deletion of the SD (security descriptor) value beneath HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree\\<TaskName> \u2014 the signature persistence-hiding trick used by the Tarrask malware, attributed by Microsoft to HAFNIUM / Silk Typhoon (G0125). With this value gone the scheduled task continues to execute but disappears from schtasks.exe and the Task Scheduler GUI.\nreferences:\n  - https://www.microsoft.com/en-us/security/blog/2022/04/12/tarrask-malware-uses-scheduled-tasks-for-defense-evasion/\n  - https://attack.mitre.org/groups/G0125/\n  - https://attack.mitre.org/software/S0987/\nauthor: SOC Detection Engineering\ndate: 2026/05/10\nlogsource:\n  category: registry_delete\n  product: windows\ndetection:\n  selection:\n    TargetObject|contains: '\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree\\'\n    TargetObject|endswith: '\\SD'\n    EventType: 'DeleteValue'\n  filter_servicing:\n    Image|endswith:\n      - '\\TrustedInstaller.exe'\n      - '\\msiexec.exe'\n  condition: selection and not filter_servicing\nfields:\n  - ComputerName\n  - User\n  - Image\n  - TargetObject\n  - EventType\nfalsepositives:\n  - Extremely rare \u2014 Windows itself does not delete the SD value for a scheduled task during normal lifecycle operations; legitimate task removal deletes the entire Tree subkey, not just the SD value.\nlevel: critical\ntags:\n  - attack.persistence\n  - attack.defense_evasion\n  - attack.t1053.005\n  - attack.t1562.001\n  - attack.t1070.009\n  - attack.g0125\n  - attack.s0987",
      "rationale": "The Tarrask malware (S0987) is, per Microsoft Threat Intelligence Center's April 2022 blog, attributed specifically to HAFNIUM / Silk Typhoon (G0125). Tarrask's distinguishing tradecraft is not creating a scheduled task \u2014 many actors do that \u2014 but selectively deleting the SD value under the task's TaskCache\\Tree key, which exploits a quirk where the Task Scheduler API treats a missing SD as 'no pe"
    }
  ]
}