# intel/.llm_actor_uc_cache/80/8037dbad0a29de09b081fb8cc450e6ffc98bcff5.json
{
  "name": "Machete",
  "ucs": [
    {
      "title": "[LLM] Machete (APT-C-43) Pyark/py2exe Python backdoor execution from user-writable paths",
      "description": "Hunts Machete's signature py2exe-compiled Python implants (Machete/Pyark/Numpy) which leave behind interpreter artifacts (python27.dll, library.zip, _ctypes.pyd) executing from %APPDATA%, %TEMP%, or %PROGRAMDATA% \u2014 frequently chained with a scheduled task for persistence. This is highly specific to Machete because almost no other tracked actor consistently ships py2exe Python loaders against LATAM military/government targets.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1059.006",
        "T1053.005",
        "T1036.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as process values(Processes.parent_process) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where (Processes.process_name IN (\"python.exe\",\"pythonw.exe\") OR Processes.process IN (\"*python27.dll*\",\"*library.zip*\",\"*_ctypes.pyd*\",\"*MSVCR90.dll*\")) AND (Processes.process_path IN (\"*\\\\AppData\\\\*\",\"*\\\\Temp\\\\*\",\"*\\\\ProgramData\\\\*\",\"*\\\\Users\\\\Public\\\\*\")) by Processes.dest Processes.process_name Processes.process Processes.parent_process_name Processes.user | `drop_dm_object_name(Processes)` | join type=inner dest [| tstats summariesonly=true count from datamodel=Endpoint.Processes where Processes.process_name=schtasks.exe AND (Processes.process=\"*python*\" OR Processes.process=\"*\\\\AppData\\\\*\" OR Processes.process=\"*\\\\Temp\\\\*\") by Processes.dest Processes.process | `drop_dm_object_name(Processes)` | rename process as schtasks_cmd] | convert ctime(firstTime) ctime(lastTime)",
      "kql": "let py_exec = DeviceProcessEvents | where Timestamp > ago(14d) | where (FileName in~ (\"python.exe\",\"pythonw.exe\")) or ProcessCommandLine has_any (\"python27.dll\",\"library.zip\",\"_ctypes.pyd\") | where FolderPath has_any (@\"\\AppData\\\", @\"\\Temp\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\") | project PyTime=Timestamp, DeviceId, DeviceName, PyProc=FileName, PyPath=FolderPath, PyCmd=ProcessCommandLine, AccountName; let sched = DeviceProcessEvents | where Timestamp > ago(14d) | where FileName =~ \"schtasks.exe\" and ProcessCommandLine has_any (\"python\",\"AppData\",\"Temp\",\"ProgramData\") and ProcessCommandLine has \"/create\" | project SchedTime=Timestamp, DeviceId, SchedCmd=ProcessCommandLine; py_exec | join kind=inner sched on DeviceId | where abs(datetime_diff('minute', PyTime, SchedTime)) <= 60 | project PyTime, SchedTime, DeviceName, AccountName, PyPath, PyCmd, SchedCmd",
      "rationale": "Machete is one of the very few espionage actors whose primary RAT (Pyark/Machete) is a py2exe-bundled Python program \u2014 leaving telltale interpreter artifacts on disk. Combining the artifact path with a near-time schtasks /create on the same host narrows to Machete's documented install\u2192persistence pattern targeting Venezuelan/Ecuadorian/Colombian military and government endpoints, sharply different"
    },
    {
      "title": "[LLM] Machete Spanish-language military-themed spearphish \u2192 RAR/ZIP self-extractor \u2192 MSIexec install chain",
      "description": "Hunts Machete's characteristic delivery chain: Spanish-language lures impersonating LATAM military/intelligence bodies (Comando, Ej\u00e9rcito, Fuerza Armada, FANB, comunicado, informe) carrying RAR/ZIP self-extractors, followed within minutes by MSIexec or a freshly-dropped binary from %TEMP%/%APPDATA%. This pivots on Machete's unique combination of Spanish military lure language and T1218.007 install behavior on LATAM tenants.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1218.007"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime from datamodel=Email where (All_Email.subject IN (\"*ejercito*\",\"*ej\u00e9rcito*\",\"*fuerza armada*\",\"*FANB*\",\"*comando*\",\"*comunicado*\",\"*informe*\",\"*orden de operaciones*\",\"*inteligencia militar*\",\"*ministerio de defensa*\")) AND (All_Email.file_name IN (\"*.rar\",\"*.zip\",\"*.lnk\",\"*.7z\",\"*.iso\")) by All_Email.src_user All_Email.recipient All_Email.subject All_Email.file_name All_Email.message_id | `drop_dm_object_name(All_Email)` | rename recipient as user | join type=inner user [| tstats summariesonly=true count values(Processes.process) as proc_chain min(_time) as proc_time from datamodel=Endpoint.Processes where (Processes.process_name=msiexec.exe AND (Processes.process=\"*http*\" OR Processes.process=\"*\\\\AppData\\\\*\" OR Processes.process=\"*\\\\Temp\\\\*\")) OR (Processes.parent_process_name IN (\"winrar.exe\",\"7zg.exe\",\"explorer.exe\") AND Processes.process_path IN (\"*\\\\AppData\\\\*\",\"*\\\\Temp\\\\*\")) by Processes.user Processes.dest Processes.process_name Processes.process | `drop_dm_object_name(Processes)`] | where proc_time >= firstTime AND proc_time <= firstTime+3600 | table firstTime proc_time user dest src_user subject file_name process_name proc_chain",
      "kql": "let lure = EmailEvents | where Timestamp > ago(7d) | where Subject matches regex @\"(?i)(ejercito|ej[e\u00e9]rcito|fuerza armada|FANB|comando|comunicado|informe|orden de operaciones|inteligencia militar|ministerio de defensa)\" | join kind=inner (EmailAttachmentInfo | where FileName matches regex @\"(?i)\\.(rar|zip|7z|iso|lnk|img)$\") on NetworkMessageId | project MailTime=Timestamp, RecipientEmailAddress, SenderFromAddress, Subject, FileName, SHA256; let sus_install = DeviceProcessEvents | where Timestamp > ago(7d) | where (FileName =~ \"msiexec.exe\" and (ProcessCommandLine has_any (\"http://\",\"https://\",@\"\\AppData\\\",@\"\\Temp\\\"))) or (InitiatingProcessFileName in~ (\"winrar.exe\",\"7zg.exe\",\"explorer.exe\") and FolderPath has_any (@\"\\AppData\\\",@\"\\Temp\\\") and FileName endswith \".exe\") | project ProcTime=Timestamp, DeviceName, AccountUpn, ProcCmd=ProcessCommandLine, ProcFile=FileName, ProcParent=InitiatingProcessFileName; lure | join kind=inner sus_install on $left.RecipientEmailAddress == $right.AccountUpn | where datetime_diff('minute', ProcTime, MailTime) between (0 .. 60) | project MailTime, ProcTime, DeviceName, AccountUpn, SenderFromAddress, Subject, FileName, ProcParent, ProcFile, ProcCmd",
      "rationale": "Machete (aka APT-C-43, El Machete) is uniquely Spanish-speaking and has consistently used military/intelligence lure language in Spanish targeting Venezuelan FANB and other LATAM defense ministries \u2014 per ESET, Cylance and Securelist reporting. Most other state actors targeting these victims (e.g. Blind Eagle/APT-C-36) prefer banking or judicial lures, so the combination of Spanish military lexicon"
    }
  ]
}