# intel/.llm_actor_uc_cache/80/8037dbad0a29de09b081fb8cc450e6ffc98bcff5.json
{
  "name": "Machete",
  "ucs": [
    {
      "title": "[LLM] Machete (APT-C-43) PyInstaller-packed RAT extraction in %TEMP%\\_MEI* (Pyark/Machete loader)",
      "description": "Hunts the defining artefact of Machete's espionage toolkit: a PyInstaller-bundled Python RAT (historically 'Pyark', 'Machete', 'Numpyrat') that, on execution, self-extracts its embedded Python runtime (python3X.dll, _ctypes.pyd, _socket.pyd, etc.) into a randomly numbered \\Temp\\_MEIxxxxx\\ scratch directory. Machete has used PyInstaller exclusively as its packer since 2014 to compile the Spanish-language Python operator code into a single Windows EXE delivered to Latin American government/military targets \u2014 the _MEI extraction pattern, paired with a non-Microsoft initiating binary running from ",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1059.006",
        "T1027.002",
        "T1204.002",
        "T1036.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime values(Filesystem.file_name) as dropped_files dc(Filesystem.file_name) as unique_drops from datamodel=Endpoint.Filesystem where Filesystem.action=created Filesystem.file_path=\"*\\\\Temp\\\\_MEI*\" by Filesystem.dest Filesystem.process_guid Filesystem.user\n| `drop_dm_object_name(Filesystem)`\n| where unique_drops >= 5\n| rename process_guid as parent_proc_guid\n| join type=inner dest parent_proc_guid [\n    | tstats `summariesonly` values(Processes.process_name) as parent_name values(Processes.process_path) as parent_path values(Processes.process) as parent_cmd values(Processes.parent_process_name) as grandparent values(Processes.process_hash) as parent_sha256 from datamodel=Endpoint.Processes by Processes.dest Processes.process_guid\n    | `drop_dm_object_name(Processes)`\n    | rename process_guid as parent_proc_guid\n  ]\n| where NOT match(parent_path, \"(?i)C:\\\\\\\\(Windows|Program Files|Program Files \\\\(x86\\\\))\\\\\\\\\")\n| where match(parent_path, \"(?i)\\\\\\\\(Users|ProgramData|AppData|Public|Temp|Downloads)\\\\\\\\\")\n| eval lure_chain=if(match(grandparent,\"(?i)(winword|excel|powerpnt|outlook|wscript|cscript|mshta)\\.exe\"),\"yes\",\"no\")\n| table firstTime dest user parent_name parent_path parent_cmd grandparent parent_sha256 unique_drops dropped_files lure_chain\n| `security_content_ctime(firstTime)`\n| sort - firstTime",
      "kql": "// Machete (APT-C-43) \u2014 PyInstaller _MEI extraction by user-context loader\nlet PyInstallerHallmarks = dynamic([\"python27.dll\",\"python36.dll\",\"python37.dll\",\"python38.dll\",\"python39.dll\",\"python310.dll\",\"python311.dll\",\"_ctypes.pyd\",\"_socket.pyd\",\"_ssl.pyd\",\"select.pyd\",\"unicodedata.pyd\",\"pyimod02_archive.pyc\",\"pyimod03_importers.pyc\",\"pyiboot01_bootstrap.pyc\"]);\nlet OfficeLures = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"acrord32.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"explorer.exe\"]);\nlet MeiDrops = DeviceFileEvents\n    | where Timestamp > ago(14d)\n    | where ActionType == \"FileCreated\"\n    | where FolderPath matches regex @\"(?i)\\\\Temp\\\\_MEI\\d{3,}\\\\\"\n    | where FileName in~ (PyInstallerHallmarks)\n          or FileName endswith \".pyd\"\n          or (FileName startswith_cs \"python\" and FileName endswith \".dll\");\nMeiDrops\n| summarize FirstDrop = min(Timestamp),\n            LastDrop  = max(Timestamp),\n            DroppedFiles = make_set(FileName, 50),\n            DistinctDropped = dcount(FileName)\n            by DeviceId, DeviceName,\n               InitiatingProcessSHA256, InitiatingProcessFileName,\n               InitiatingProcessFolderPath, InitiatingProcessCommandLine,\n               InitiatingProcessAccountName, InitiatingProcessParentFileName\n| where DistinctDropped >= 3                                  // PyInstaller drops 5-50 files; require >=3 to dodge legit single-DLL drops\n| where InitiatingProcessAccountName !endswith \"$\"            // exclude SYSTEM / machine-account installers\n| where not(InitiatingProcessFolderPath startswith_cs @\"C:\\Windows\\\")\n      and not(InitiatingProcessFolderPath startswith_cs @\"C:\\Program Files\\\")\n      and not(InitiatingProcessFolderPath startswith_cs @\"C:\\Program Files (x86)\\\")\n| where InitiatingProcessFolderPath has_any (@\"\\Users\\\", @\"\\AppData\\\", @\"\\ProgramData\\\", @\"\\Public\\\", @\"\\Temp\\\", @\"\\Downloads\\\")\n| extend OfficeLureChain = iif(tolower(InitiatingProcessParentFileName) in (OfficeLures), \"yes\", \"no\")\n| project FirstDrop, LastDrop, DeviceName, InitiatingProcessAccountName,\n          LoaderBinary = InitiatingProcessFileName,\n          LoaderPath   = InitiatingProcessFolderPath,\n          LoaderCmd    = InitiatingProcessCommandLine,\n          LoaderSHA256 = InitiatingProcessSHA256,\n          Grandparent  = InitiatingProcessParentFileName,\n          OfficeLureChain, DistinctDropped, DroppedFiles\n| order by FirstDrop desc",
      "sentinel_kql": "// Machete (APT-C-43) \u2014 PyInstaller _MEI extraction (ASIM port)\nlet PyInstallerHallmarks = dynamic([\"python27.dll\",\"python36.dll\",\"python37.dll\",\"python38.dll\",\"python39.dll\",\"python310.dll\",\"python311.dll\",\"_ctypes.pyd\",\"_socket.pyd\",\"_ssl.pyd\",\"select.pyd\",\"unicodedata.pyd\",\"pyimod02_archive.pyc\",\"pyimod03_importers.pyc\",\"pyiboot01_bootstrap.pyc\"]);\nImFileEvent\n| where TimeGenerated > ago(14d)\n| where EventType == \"FileCreated\"\n| where TargetFilePath matches regex @\"(?i)\\\\Temp\\\\_MEI\\d{3,}\\\\\"\n| where TargetFileName in~ (PyInstallerHallmarks)\n      or TargetFileName endswith \".pyd\"\n      or (TargetFileName startswith_cs \"python\" and TargetFileName endswith \".dll\")\n| where ActorUsername !endswith \"$\"\n| summarize FirstDrop = min(TimeGenerated),\n            LastDrop  = max(TimeGenerated),\n            DroppedFiles = make_set(TargetFileName, 50),\n            DistinctDropped = dcount(TargetFileName)\n            by DvcHostname, ActorUsername,\n               ActingProcessName, ActingProcessCommandLine\n| where DistinctDropped >= 3\n| where not(ActingProcessName has_any (@\"C:\\Windows\\\", @\"C:\\Program Files\\\"))\n| order by FirstDrop desc",
      "sigma_yaml": "title: Machete APT-C-43 PyInstaller _MEI Runtime Extraction\nid: 5e6f3b27-4a9c-46d8-a4d5-8c1ae5b1d9ce\nstatus: experimental\ndescription: Detects creation of PyInstaller bootstrap files (python3X.dll, _ctypes.pyd, pyimod*.pyc) under %TEMP%\\_MEIxxxxx\\ \u2014 the unpacking pattern of Machete (APT-C-43) Spanish-language Python RATs (Pyark / Machete / Numpyrat).\nreferences:\n  - https://attack.mitre.org/groups/G0095/\n  - https://www.welivesecurity.com/2019/08/05/machete-just-go-through/\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.execution\n  - attack.t1059.006\n  - attack.defense_evasion\n  - attack.t1027.002\n  - attack.g0095\nlogsource:\n  category: file_event\n  product: windows\ndetection:\n  selection_path:\n    TargetFilename|re: '(?i)\\\\Temp\\\\_MEI\\d{3,}\\\\'\n  selection_artefact:\n    TargetFilename|endswith:\n      - '\\python27.dll'\n      - '\\python36.dll'\n      - '\\python37.dll'\n      - '\\python38.dll'\n      - '\\python39.dll'\n      - '\\python310.dll'\n      - '\\python311.dll'\n      - '\\_ctypes.pyd'\n      - '\\_socket.pyd'\n      - '\\_ssl.pyd'\n      - '\\pyimod02_archive.pyc'\n      - '\\pyimod03_importers.pyc'\n      - '\\pyiboot01_bootstrap.pyc'\n  filter_loader_path:\n    Image|startswith:\n      - 'C:\\Windows\\'\n      - 'C:\\Program Files\\'\n      - 'C:\\Program Files (x86)\\'\n  condition: selection_path and selection_artefact and not filter_loader_path\nfalsepositives:\n  - Legitimate PyInstaller-packed enterprise tooling (Splunk UF wrappers, vendor agents) running from user paths \u2014 baseline by SHA256 of the loader\nlevel: high",
      "rationale": "Machete is one of very few state-aligned actors whose entire offensive kit (publicly tracked since 2014 by Kaspersky, ESET, Cylance) is built in Python and shipped as PyInstaller-packed EXEs \u2014 when the loader fires, PyInstaller deterministically self-extracts a Python interpreter and stdlib .pyd modules into %TEMP%\\_MEIxxxxx\\. Generic 'suspicious python' rules miss this because Machete operators n"
    },
    {
      "title": "[LLM] Machete LatAm spearphish: Spanish-lure Office attachment \u2192 scripting host \u2192 schtasks persistence",
      "description": "Hunts Machete's end-to-end intrusion pattern observed against Venezuelan/Ecuadorian/Colombian government and military targets: an inbound Spanish-language email with an Office or RAR/SFX attachment is opened, the Office host (winword/excel/powerpnt) spawns a scripting interpreter (cscript/wscript/cmd/python), and within minutes a scheduled task is registered via schtasks.exe \u2014 Machete's documented persistence mechanism (T1053.005) for the Pyark/Machete RAT. The Spanish-keyword anchor in the lure plus the Office\u2192script\u2192schtasks chain on the same host is highly specific to APT-C-43.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.003",
        "T1059.005",
        "T1053.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` count min(_time) as schtasks_time values(Processes.process) as schtasks_cmd values(Processes.user) as user values(Processes.parent_process_name) as scripting_host from datamodel=Endpoint.Processes where Processes.process_name=schtasks.exe (Processes.process=\"*/create*\" OR Processes.process=\"*/CREATE*\") Processes.parent_process_name IN (cmd.exe,wscript.exe,cscript.exe,powershell.exe,python.exe,pythonw.exe,mshta.exe) by Processes.dest Processes.parent_process_guid\n| `drop_dm_object_name(Processes)`\n| rename parent_process_guid as scripting_guid\n| join type=inner dest scripting_guid [\n    | tstats `summariesonly` min(_time) as script_time values(Processes.parent_process_name) as office_host values(Processes.process) as script_cmd from datamodel=Endpoint.Processes where Processes.parent_process_name IN (winword.exe,excel.exe,powerpnt.exe,outlook.exe,acrord32.exe,wscript.exe,wmplayer.exe) by Processes.dest Processes.process_guid\n    | `drop_dm_object_name(Processes)`\n    | rename process_guid as scripting_guid\n  ]\n| where (schtasks_time - script_time) < 600 AND schtasks_time >= script_time\n| join type=left user [\n    | tstats `summariesonly` values(All_Email.subject) as subject values(All_Email.src_user) as sender from datamodel=Email.All_Email where All_Email.action=delivered (All_Email.subject=\"*factura*\" OR All_Email.subject=\"*comunicado*\" OR All_Email.subject=\"*ministerio*\" OR All_Email.subject=\"*ejercito*\" OR All_Email.subject=\"*ej\\xe9rcito*\" OR All_Email.subject=\"*fuerza armada*\" OR All_Email.subject=\"*denuncia*\" OR All_Email.subject=\"*urgente*\" OR All_Email.subject=\"*confidencial*\" OR All_Email.subject=\"*citaci\\xf3n*\") by All_Email.recipient\n    | `drop_dm_object_name(All_Email)`\n    | rename recipient as user\n  ]\n| table schtasks_time dest user office_host scripting_host script_cmd schtasks_cmd subject sender\n| sort - schtasks_time",
      "kql": "// Machete (APT-C-43) \u2014 Spanish-lure phish \u2192 Office \u2192 script host \u2192 schtasks persistence chain\nlet WindowMin = 10m;\nlet OfficeLures   = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"acrord32.exe\"]);\nlet ScriptHosts   = dynamic([\"cmd.exe\",\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\",\"python.exe\",\"pythonw.exe\"]);\nlet SpanishLureKw = dynamic([\"factura\",\"comunicado\",\"ministerio\",\"ejercito\",\"ej\u00e9rcito\",\"fuerza armada\",\"denuncia\",\"urgente\",\"confidencial\",\"citaci\u00f3n\",\"citacion\",\"presidencia\",\"asamblea\",\"informe\",\"oficio\"]);\n// Stage 1: Spanish-language inbound phishing email with attachment\nlet PhishMail = EmailEvents\n    | where Timestamp > ago(7d)\n    | where EmailDirection == \"Inbound\"\n          and DeliveryAction == \"Delivered\"\n          and AttachmentCount >= 1\n    | where (EmailLanguage in~ (\"es\",\"es-ES\",\"es-MX\",\"es-VE\",\"es-CO\",\"es-AR\",\"pt-BR\")\n            or Subject has_any (SpanishLureKw))\n    | project EmailTime = Timestamp, NetworkMessageId, Subject,\n              SenderFromAddress, SenderFromDomain, RecipientEmailAddress;\n// Stage 2: schtasks /create launched by a scripting host on the recipient's machine\nlet PersistChain = DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where FileName =~ \"schtasks.exe\"\n    | where ProcessCommandLine has_any (\"/create\",\"/Create\",\"/CREATE\")\n    | where AccountName !endswith \"$\"\n    | where InitiatingProcessFileName in~ (ScriptHosts)\n    | join kind=inner (\n        DeviceProcessEvents\n        | where Timestamp > ago(7d)\n        | where InitiatingProcessFileName in~ (OfficeLures)\n        | where FileName in~ (ScriptHosts)\n        | project ScriptTime = Timestamp, DeviceId,\n                  ScriptHost = FileName, ScriptCmd = ProcessCommandLine,\n                  ScriptProcId = ProcessId, ScriptCreationTime = ProcessCreationTime,\n                  OfficeHost = InitiatingProcessFileName,\n                  OfficeCmd  = InitiatingProcessCommandLine,\n                  Recipient  = AccountUpn\n      ) on DeviceId\n              and $left.InitiatingProcessId == $right.ScriptProcId\n              and $left.InitiatingProcessCreationTime == $right.ScriptCreationTime\n    | where Timestamp between (ScriptTime .. ScriptTime + WindowMin)\n    | project SchtasksTime = Timestamp, DeviceName, AccountName, AccountUpn,\n              OfficeHost, OfficeCmd, ScriptHost, ScriptCmd,\n              SchtasksCmd = ProcessCommandLine, Recipient;\n// Stage 3: correlate to a Spanish-language phish to the same recipient within 24h prior\nPersistChain\n| join kind=leftouter PhishMail on $left.Recipient == $right.RecipientEmailAddress\n| where isnull(EmailTime) or (SchtasksTime between (EmailTime .. EmailTime + 24h))\n| project SchtasksTime, DeviceName, AccountUpn,\n          OfficeHost, OfficeCmd, ScriptHost, ScriptCmd, SchtasksCmd,\n          PhishSubject = Subject, PhishSender = SenderFromAddress, PhishTime = EmailTime\n| order by SchtasksTime desc",
      "sentinel_kql": "// Machete (APT-C-43) \u2014 Office \u2192 script host \u2192 schtasks chain (Sentinel SecurityEvent + OfficeActivity)\nlet WindowMin = 10m;\nlet OfficeLures   = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"acrord32.exe\"]);\nlet ScriptHosts   = dynamic([\"cmd.exe\",\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\",\"python.exe\",\"pythonw.exe\"]);\nlet SpanishLureKw = dynamic([\"factura\",\"comunicado\",\"ministerio\",\"ejercito\",\"ej\u00e9rcito\",\"fuerza armada\",\"denuncia\",\"urgente\",\"confidencial\",\"citaci\u00f3n\",\"citacion\",\"presidencia\",\"asamblea\",\"informe\",\"oficio\"]);\nlet ScriptFromOffice = SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ChildExe = tolower(tostring(split(NewProcessName,\"\\\\\")[-1])),\n             ParentExe = tolower(tostring(split(ParentProcessName,\"\\\\\")[-1]))\n    | where ParentExe in (OfficeLures) and ChildExe in (ScriptHosts)\n    | project ScriptTime = TimeGenerated, Computer, Account,\n              OfficeHost = ParentExe, ScriptHost = ChildExe,\n              ScriptPid = NewProcessId, ScriptCmd = CommandLine;\nlet SchtasksByScript = SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ChildExe = tolower(tostring(split(NewProcessName,\"\\\\\")[-1])),\n             ParentExe = tolower(tostring(split(ParentProcessName,\"\\\\\")[-1]))\n    | where ChildExe == \"schtasks.exe\"\n    | where CommandLine has_any (\"/create\",\"/Create\",\"/CREATE\")\n    | where ParentExe in (ScriptHosts)\n    | project SchtasksTime = TimeGenerated, Computer, Account,\n              ParentPid = ProcessId, SchtasksCmd = CommandLine;\nlet Chain = ScriptFromOffice\n    | join kind=inner SchtasksByScript on Computer, $left.ScriptPid == $right.ParentPid\n    | where SchtasksTime between (ScriptTime .. ScriptTime + WindowMin);\nlet PhishMail = OfficeActivity\n    | where TimeGenerated > ago(7d)\n    | where OfficeWorkload =~ \"Exchange\"\n    | where Operation in~ (\"Deliver\",\"MailItemsAccessed\")\n    | where Subject has_any (SpanishLureKw)\n    | project PhishTime = TimeGenerated, Recipient = MailboxOwnerUPN,\n              Subject, Sender;\nChain\n| extend AccountUpn = tolower(strcat(tostring(split(Account,\"\\\\\")[-1]),\"@\",tostring(split(Account,\"\\\\\")[0])))\n| join kind=leftouter PhishMail on $left.Account == $right.Recipient\n| where isnull(PhishTime) or (SchtasksTime between (PhishTime .. PhishTime + 24h))\n| project SchtasksTime, Computer, Account, OfficeHost, ScriptHost, ScriptCmd, SchtasksCmd, PhishSubject = Subject, PhishSender = Sender, PhishTime\n| order by SchtasksTime desc",
      "sigma_yaml": "",
      "rationale": "ESET ('Machete just got sharper', 2019) and Cylance ('El Machete's Malware Attacks Cut Through LATAM', 2017) both document the same intrusion pattern: Spanish-language phish \u2192 Office/RAR self-extractor \u2192 cscript/cmd dropper \u2192 `schtasks /create` registering the PyInstaller-packed Pyark/Machete RAT for persistence. Each stage on its own (Office spawning script, schtasks creation) is too noisy to ale"
    }
  ]
}