# intel/.llm_actor_uc_cache/f0/f0306c9e17fdd9c5453348eef85c96c8b12ad457.json
{
  "name": "APT18",
  "ucs": [
    {
      "title": "[LLM] APT18 (Dynamite Panda) legacy at.exe scheduled-task persistence chained to cmd.exe payload",
      "description": "Hunts the deprecated `at.exe` Windows scheduling utility being used to register a job that subsequently spawns a Windows command interpreter or scripting host. APT18 / TG-0416 is documented (MITRE G0026 \u2014 T1053.002 + T1059.003) to lean on the legacy `at` command rather than `schtasks.exe` for persistence and timed execution \u2014 an idiomatic Chinese-state TTP that survives because `at.exe` is largely abandoned by legitimate admins on modern Windows.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1053.002",
        "T1059.003"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as process values(Processes.parent_process) as parent_process from datamodel=Endpoint.Processes where (Processes.process_name=\"at.exe\" AND (match(Processes.process,\"(?i)\\b\\d{1,2}:\\d{2}\\b\") OR Processes.process IN (\"*/interactive*\",\"*/every:*\",\"*/next:*\"))) OR (Processes.parent_process_name=\"at.exe\" AND Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\")) by host Processes.user Processes.process_name Processes.parent_process_name Processes.process Processes.parent_process | `drop_dm_object_name(Processes)` | where NOT match(user,\"\\$$\") | sort - lastTime",
      "kql": "let LookbackDays = 30d;\n// (1) Direct at.exe scheduling \u2014 schedule-time argument, /interactive flag, or /every:/next: switches\nlet AtUsage = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where FileName =~ \"at.exe\"\n    | where InitiatingProcessFolderPath startswith @\"C:\\Windows\\\"     // legit on-disk at.exe only\n    | where ProcessCommandLine matches regex @\"(?i)\\b\\d{1,2}:\\d{2}\\b\"\n         or ProcessCommandLine has_any (\"/interactive\",\"/every:\",\"/next:\");\n// (2) Children spawned by an at-registered job (parent = at.exe)\nlet AtChildren = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where InitiatingProcessFileName =~ \"at.exe\"\n    | where FileName in~ (\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\");\nunion AtUsage, AtChildren\n| project Timestamp, DeviceName, AccountName,\n          FileName, ProcessCommandLine,\n          ParentImage = InitiatingProcessFileName,\n          ParentCmd   = InitiatingProcessCommandLine,\n          IsRemoteSession = AdditionalFields\n| order by Timestamp desc",
      "sentinel_kql": "let LookbackDays = 30d;\nSecurityEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventID == 4688                       // process create \u2014 requires audit policy + cmdline GPO\n| where Account !endswith \"$\"\n| where (NewProcessName endswith \"\\\\at.exe\"\n          and (CommandLine matches regex @\"(?i)\\b\\d{1,2}:\\d{2}\\b\"\n               or CommandLine has_any (\"/interactive\",\"/every:\",\"/next:\")))\n     or (ParentProcessName endswith \"\\\\at.exe\"\n          and NewProcessName matches regex @\"(?i)\\\\(cmd|powershell|pwsh|cscript|wscript|mshta|rundll32|regsvr32)\\.exe$\")\n| project TimeGenerated, Computer, Account,\n          NewProcessName, CommandLine,\n          ParentProcessName, TokenElevationType, LogonId\n| order by TimeGenerated desc",
      "sigma_yaml": "title: APT18 (Dynamite Panda) Legacy at.exe Scheduled Task Abuse\nid: 9c3d2f5b-1a44-4f1c-b6c3-1a2e7d4f8b21\nstatus: experimental\ndescription: Detects use of the deprecated `at.exe` scheduling utility (or processes spawned by it) \u2014 a tradecraft signature of APT18 / TG-0416 / Dynamite Panda for timed execution and persistence on Windows hosts.\nreferences:\n  - https://attack.mitre.org/groups/G0026/\n  - https://attack.mitre.org/techniques/T1053/002/\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.persistence\n  - attack.execution\n  - attack.t1053.002\n  - attack.t1059.003\n  - attack.g0026\nlogsource:\n  product: windows\n  category: process_creation\ndetection:\n  selection_at_schedule:\n    Image|endswith: '\\at.exe'\n  selection_at_args:\n    - CommandLine|re: '(?i)\\b\\d{1,2}:\\d{2}\\b'\n    - CommandLine|contains:\n        - '/interactive'\n        - '/every:'\n        - '/next:'\n  selection_at_child:\n    ParentImage|endswith: '\\at.exe'\n    Image|endswith:\n      - '\\cmd.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\cscript.exe'\n      - '\\wscript.exe'\n      - '\\mshta.exe'\n      - '\\rundll32.exe'\n      - '\\regsvr32.exe'\n  filter_machine_account:\n    User|endswith: '$'\n  condition: ((selection_at_schedule and selection_at_args) or selection_at_child) and not filter_machine_account\nfalsepositives:\n  - Legacy administrator scripts on pre-Windows-10 hosts\n  - Backup or monitoring software still scripted around at.exe\nlevel: high",
      "rationale": "APT18 is documented in MITRE G0026 as one of the few well-tracked clusters that explicitly uses T1053.002 (`at` command) over T1053.005 (`schtasks`). On modern Windows, `at.exe` was deprecated in Windows 8 / Server 2012 and is essentially absent from legitimate admin workflows \u2014 so `at.exe` invoked with a HH:MM schedule, /interactive, /every: or /next:, OR `at.exe` parenting cmd.exe / scripting ho"
    },
    {
      "title": "[LLM] APT18 dropper-in-user-path \u2192 Run-key persistence \u2192 external HTTP beacon (HTTPBrowser/Gh0st chain)",
      "description": "Correlates three APT18 staples within a one-hour window on the same host: (a) a binary executing from a user-writable path (\\AppData\\, \\ProgramData\\, \\Windows\\Temp\\, \\Public\\), (b) that same process \u2014 or a sibling under the same path \u2014 writing a Run/RunOnce value, and (c) outbound traffic to a public IP. APT18's HTTPBrowser and Gh0st-derived RATs (T1547.001 + T1071.001 + T1105) historically drop into user-writable directories, register Run-key autostart, and beacon over HTTP/HTTPS to commodity-hosted C2.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "medium",
      "techs": [
        "T1547.001",
        "T1071.001",
        "T1105",
        "T1027.013"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as persistTime from datamodel=Endpoint.Registry where Registry.registry_path IN (\"*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run*\",\"*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce*\") AND match(Registry.registry_value_data,\"(?i)(\\\\Users\\\\[^\\\\\\\\]+\\\\AppData\\\\|\\\\ProgramData\\\\|\\\\Windows\\\\Temp\\\\|\\\\Public\\\\)\") AND Registry.process_name!=\"msiexec.exe\" AND Registry.process_name!=\"explorer.exe\" by host Registry.user Registry.registry_path Registry.registry_value_name Registry.registry_value_data Registry.process_name | `drop_dm_object_name(Registry)` | where NOT match(user,\"\\$$\") | join type=inner host [| tstats summariesonly=true count values(All_Traffic.dest_ip) as dest_ip values(All_Traffic.dest_port) as dest_port values(All_Traffic.app) as c2_url from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest_category!=\"private\" AND match(All_Traffic.process_path,\"(?i)\\\\Users\\\\[^\\\\\\\\]+\\\\AppData\\\\|\\\\ProgramData\\\\|\\\\Windows\\\\Temp\\\\|\\\\Public\\\\\") by host All_Traffic.process_name All_Traffic.process_path _time | `drop_dm_object_name(All_Traffic)` | rename _time as netTime] | where abs(netTime - persistTime) <= 3600 | table host user registry_path registry_value_data process_name dest_ip dest_port c2_url persistTime netTime",
      "kql": "let LookbackDays = 14d;\nlet WindowMin = 60;\nlet UserWritable = @\"(?i)\\\\Users\\\\[^\\\\]+\\\\AppData\\\\|\\\\ProgramData\\\\|\\\\Windows\\\\Temp\\\\|\\\\Public\\\\\";\n// (1) Run/RunOnce value pointing into a user-writable path, set by a non-installer process\nlet RunKeyDrops = DeviceRegistryEvents\n    | where Timestamp > ago(LookbackDays)\n    | where ActionType == \"RegistryValueSet\"\n    | where RegistryKey has_any (@\"\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\",\n                                  @\"\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\")\n    | where InitiatingProcessAccountName !endswith \"$\"\n    | where InitiatingProcessFileName !in~ (\"msiexec.exe\",\"explorer.exe\",\"setup.exe\",\"installer.exe\",\"trustedinstaller.exe\")\n    | where RegistryValueData matches regex UserWritable\n    | project PersistTime = Timestamp, DeviceId, DeviceName,\n              PersistUser = InitiatingProcessAccountName,\n              PersistParent = InitiatingProcessFileName,\n              RegistryKey, RegistryValueName, RegistryValueData;\n// (2) Same host: process running from user-writable path beaconing to public IP\nDeviceNetworkEvents\n| where Timestamp > ago(LookbackDays)\n| where RemoteIPType == \"Public\"\n| where InitiatingProcessFolderPath matches regex UserWritable\n| where InitiatingProcessAccountName !endswith \"$\"\n| join kind=inner RunKeyDrops on DeviceId\n| where Timestamp between (PersistTime - WindowMin*1m .. PersistTime + WindowMin*1m)\n| project Timestamp, DeviceName, PersistUser,\n          BeaconProcess = InitiatingProcessFileName,\n          BeaconPath    = InitiatingProcessFolderPath,\n          BeaconCmd     = InitiatingProcessCommandLine,\n          RemoteIP, RemotePort, RemoteUrl,\n          RegistryKey, RegistryValueName, RegistryValueData,\n          PersistTime,\n          DelayMin = datetime_diff('minute', Timestamp, PersistTime)\n| order by Timestamp desc",
      "sentinel_kql": "let LookbackDays = 14d;\nlet WindowMin = 60;\nlet UserWritable = @\"(?i)\\\\Users\\\\[^\\\\]+\\\\AppData\\\\|\\\\ProgramData\\\\|\\\\Windows\\\\Temp\\\\|\\\\Public\\\\\";\nlet RunKeyDrops = ImRegistryEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n    | where RegistryKey has_any (@\"\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\",\n                                  @\"\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\")\n    | where ActorUsername !endswith \"$\"\n    | where ActingProcessName !in~ (\"msiexec.exe\",\"explorer.exe\",\"setup.exe\",\"installer.exe\",\"trustedinstaller.exe\")\n    | where RegistryValueData matches regex UserWritable\n    | project PersistTime = TimeGenerated, DvcHostname, ActorUsername,\n              ActingProcessName, RegistryKey, RegistryValue, RegistryValueData;\nImNetworkSession\n| where TimeGenerated > ago(LookbackDays)\n| where NetworkDirection =~ \"Outbound\" or isempty(NetworkDirection)\n| where not(ipv4_is_private(DstIpAddr))\n| where ActingProcessCommandLine matches regex UserWritable\n   or ActingProcessName matches regex UserWritable\n| join kind=inner RunKeyDrops on $left.DvcHostname == $right.DvcHostname\n| where TimeGenerated between (PersistTime - WindowMin*1m .. PersistTime + WindowMin*1m)\n| project TimeGenerated, DvcHostname, ActorUsername,\n          BeaconProcess = ActingProcessName,\n          BeaconCmd     = ActingProcessCommandLine,\n          DstIpAddr, DstPortNumber, Url,\n          RegistryKey, RegistryValue, RegistryValueData, PersistTime,\n          DelayMin = datetime_diff('minute', TimeGenerated, PersistTime)\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "This is the canonical APT18/Dynamite Panda chain documented in CrowdStrike and Dell SecureWorks reporting on TG-0416: drop a Gh0st/HTTPBrowser-derived implant into AppData/ProgramData/Temp, hook it into HKCU/HKLM Run for autostart (T1547.001), then HTTP-beacon to public infrastructure (T1071.001 + T1105). Each pillar individually is noisy (any installer touches Run keys, any service connects to pu"
    }
  ]
}