# intel/.llm_actor_uc_cache/f0/f0bbfabf1ec6230dc85284983c871814a24e96e0.json
{
  "name": "RansomHub",
  "ucs": [
    {
      "title": "[LLM] RansomHub EDRKillShifter BYOVD driver-load to terminate EDR",
      "description": "Hunts for the EDRKillShifter loader and its signed-but-vulnerable drivers (e.g., RentDrv2.sys, TrueSight.sys, ThreatFireMonitor.sys) that RansomHub affiliates run pre-encryption to kill EDR/AV. The combo of an unsigned user-mode loader dropping a signed kernel driver into a non-standard path, then registering it as a service, is RansomHub-characteristic and was specifically called out by Sophos X-Ops, CrowdStrike, and CISA AA24-242A.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1562.001",
        "T1068",
        "T1543.003",
        "T1027.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.process_path) as process_path from datamodel=Endpoint.Processes where Processes.process_name IN (\"EDRKillShifter.exe\",\"loader.exe\",\"killer.exe\",\"shifter.exe\") OR (Processes.process=\"*config.txt*\" AND Processes.process=\"*pass*\") OR Processes.process IN (\"*RentDrv2.sys*\",\"*TrueSight.sys*\",\"*ThreatFireMonitor.sys*\",\"*viragt64.sys*\",\"*WnBios.sys*\",\"*aswArPot.sys*\") by host Processes.user Processes.parent_process_name Processes.process_name Processes.process_id _time \n| `drop_dm_object_name(Processes)` \n| join type=outer host [| tstats summariesonly=t count from datamodel=Endpoint.Filesystem where (Filesystem.file_name IN (\"RentDrv2.sys\",\"TrueSight.sys\",\"ThreatFireMonitor.sys\",\"viragt64.sys\",\"WnBios.sys\",\"aswArPot.sys\") AND Filesystem.file_path!=\"*\\\\System32\\\\drivers\\\\*\") by host Filesystem.file_path Filesystem.file_name | `drop_dm_object_name(Filesystem)`] \n| join type=outer host [| tstats summariesonly=t count from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\Services\\\\*\" AND Registry.registry_value_name=\"ImagePath\" AND Registry.registry_value_data IN (\"*RentDrv*\",\"*TrueSight*\",\"*ThreatFireMonitor*\",\"*viragt64*\") by host Registry.registry_key_name | `drop_dm_object_name(Registry)`] \n| where isnotnull(file_name) OR isnotnull(registry_key_name) OR match(process_name,\"(?i)edrkillshifter\") \n| sort - lastTime",
      "kql": "let vulnDrivers = dynamic([\"rentdrv2.sys\",\"truesight.sys\",\"threatfiremonitor.sys\",\"viragt64.sys\",\"wnbios.sys\",\"aswarpot.sys\",\"dbutil_2_3.sys\"]);\nlet loaderNames = dynamic([\"edrkillshifter.exe\",\"shifter.exe\",\"killer.exe\",\"loader.exe\"]);\nlet driverDrops =\n    DeviceFileEvents\n    | where Timestamp > ago(7d)\n    | where FileName has_any (vulnDrivers)\n    | where FolderPath !has @\"\\System32\\drivers\\\" and FolderPath !has @\"\\SysWOW64\\drivers\\\"\n    | project Timestamp, DeviceId, DeviceName, DroppedDriver=FileName, DropPath=FolderPath, InitiatingProcessFileName, InitiatingProcessCommandLine, InitiatingProcessAccountName;\nlet loaderExec =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where FileName in~ (loaderNames)\n       or (ProcessCommandLine matches regex @\"(?i)\\.exe\\s+[A-Fa-f0-9]{32,}\\b\")  // hex passphrase arg pattern\n       or ProcessCommandLine has_all (\"config\",\".txt\")\n       or ProcessCommandLine has_any (vulnDrivers)\n    | project Timestamp, DeviceId, DeviceName, ProcLoader=FileName, LoaderPath=FolderPath, ProcessCommandLine, AccountName, ParentProcessFileName;\nlet svcInstall =\n    DeviceRegistryEvents\n    | where Timestamp > ago(7d)\n    | where RegistryKey has @\"\\System\\CurrentControlSet\\Services\\\"\n    | where RegistryValueName == \"ImagePath\"\n    | where RegistryValueData has_any (vulnDrivers)\n    | project Timestamp, DeviceId, DeviceName, RegistryKey, RegistryValueData;\ndriverDrops\n| join kind=inner (loaderExec) on DeviceId\n| join kind=leftouter (svcInstall) on DeviceId\n| where abs(datetime_diff('minute', Timestamp, Timestamp1)) < 30\n| project Timestamp, DeviceName, AccountName, ProcLoader, LoaderPath, ProcessCommandLine, DroppedDriver, DropPath, RegistryValueData",
      "rationale": "EDRKillShifter is a near-unique RansomHub affiliate tool \u2014 Sophos first attributed it to a RansomHub intrusion in mid-2024, and CISA AA24-242A reiterated it. The detection chains three observable artefacts that legitimate admins almost never produce together: (1) a loader binary executing with a hex passphrase argument or config.txt, (2) a known vulnerable signed driver dropped outside %SystemRoot"
    },
    {
      "title": "[LLM] RansomHub pre-encryption rclone exfil to MEGA/S3 followed by VSS deletion",
      "description": "Hunts the RansomHub data-theft-then-encrypt sequence: affiliates routinely stage data with rclone.exe (often renamed) pushing to mega.nz, backblaze, or S3, then immediately invoke vssadmin/wmic to delete shadow copies and stop services before the Linux/Windows encryptor runs. CISA AA24-242A and Mandiant track this exact rclone-to-MEGA \u2192 vssadmin delete shadows \u2192 bcdedit chain as RansomHub TTPs.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1572",
        "T1567.002",
        "T1490",
        "T1489",
        "T1486"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdlines values(Processes.process_name) as procs from datamodel=Endpoint.Processes where (Processes.process IN (\"*rclone*\",\"*\\\\svchost.exe* copy *remote:*\",\"* copy * mega:*\",\"* sync * --transfers*\",\"*config=* --no-check-certificate*\") OR Processes.process IN (\"*vssadmin* delete shadows*\",\"*wmic* shadowcopy delete*\",\"*bcdedit* /set*recoveryenabled* No*\",\"*bcdedit* /set*bootstatuspolicy* ignoreallfailures*\",\"*wevtutil* cl *\",\"*net stop*\",\"*sc* stop*\",\"*Stop-Service*\")) by host Processes.user Processes.process_name Processes.parent_process_name _time span=1h \n| `drop_dm_object_name(Processes)` \n| stats values(procs) as procs values(cmdlines) as cmdlines dc(eval(case(match(cmdlines,\"(?i)rclone|mega:|backblaze|s3:|onedrive:|--transfers\"),\"exfil\",match(cmdlines,\"(?i)vssadmin|shadowcopy|bcdedit|wevtutil\"),\"recovery\",match(cmdlines,\"(?i)net stop|sc stop|Stop-Service\"),\"svcstop\"))) as stages min(firstTime) as firstTime max(lastTime) as lastTime by host user \n| where stages>=2 \n| eval window_minutes=round((lastTime-firstTime)/60,0) \n| where window_minutes <= 240 \n| `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` \n| sort - lastTime",
      "kql": "let exfilWindow = 4h;\nlet rcloneSigns = dynamic([\"rclone\",\"--config\",\" copy \",\" sync \",\"mega:\",\"backblaze:\",\"s3:\",\"onedrive:\",\"--transfers\",\"--no-check-certificate\"]);\nlet recoverySigns = dynamic([\"vssadmin delete shadows\",\"wmic shadowcopy delete\",\"bcdedit /set\",\"recoveryenabled no\",\"bootstatuspolicy ignoreallfailures\",\"wevtutil cl\"]);\nlet svcStop = dynamic([\"net stop \",\"sc.exe stop\",\"Stop-Service\",\"taskkill /im sql\",\"taskkill /im veeam\",\"taskkill /im backup\"]);\nlet exfil =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where (FileName =~ \"rclone.exe\" or InitiatingProcessFileName =~ \"rclone.exe\")\n         or ProcessCommandLine has_any (rcloneSigns)\n    | where ProcessCommandLine has_any (dynamic([\"mega:\",\"backblaze:\",\"s3:\",\"onedrive:\",\"gdrive:\",\"dropbox:\",\"--transfers\",\"copy \",\"sync \"]))\n    | project ExfilTime=Timestamp, DeviceId, DeviceName, ExfilCmd=ProcessCommandLine, ExfilUser=AccountName, ExfilProc=FileName;\nlet recovery =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where ProcessCommandLine has_any (recoverySigns)\n    | project RecTime=Timestamp, DeviceId, RecCmd=ProcessCommandLine;\nlet stops =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where ProcessCommandLine has_any (svcStop)\n    | where ProcessCommandLine has_any (dynamic([\"sql\",\"veeam\",\"backup\",\"exchange\",\"vss\",\"sophos\",\"defender\",\"mssql\"]))\n    | project StopTime=Timestamp, DeviceId, StopCmd=ProcessCommandLine;\nexfil\n| join kind=inner (recovery) on DeviceId\n| where datetime_diff('minute', RecTime, ExfilTime) between (-240 .. 240)\n| join kind=leftouter (stops) on DeviceId\n| summarize FirstSeen=min(ExfilTime), LastSeen=max(RecTime),\n            ExfilCommands=make_set(ExfilCmd,5),\n            RecoveryCommands=make_set(RecCmd,5),\n            ServiceStopCommands=make_set(StopCmd,5),\n            User=any(ExfilUser)\n            by DeviceName, DeviceId\n| extend Stages = iff(array_length(ServiceStopCommands)>0, 3, 2)\n| where Stages >= 2\n| order by LastSeen desc",
      "rationale": "The rclone-to-cloud + shadow-copy-deletion + service-stop trio in a <4h window is the canonical RansomHub double-extortion playbook documented in CISA AA24-242A, Symantec, and Mandiant reporting. Individually each command (vssadmin delete, net stop, rclone) can be benign or noisy, but their co-occurrence on a single host within hours is essentially diagnostic of human-operated ransomware staging \u2014"
    }
  ]
}