# intel/.llm_actor_uc_cache/6f/6fe96a5b44a3c74e0f51511bf1a998f4f9d7e81f.json
{
  "name": "APT32",
  "ucs": [
    {
      "title": "[LLM] APT32 / OceanLotus DLL side-loading via signed third-party binary relocated to user-writable path",
      "description": "APT32 (OceanLotus / Cobalt Kitty / Canvas Cyclone / BISMUTH) is famous for dropping a legitimate signed binary (Google Updater, Symantec RasTls, McAfee mcoemcpy, Sysinternals/Adobe utilities, etc.) into a user-writable directory together with a trojanised DLL that the binary will side-load (KerrDown, Denis, Cobalt Strike stager). This UC catches the side-load by surfacing those vendor binaries when they execute from non-default paths such as %APPDATA%, %ProgramData%, %Public% or %TEMP% \u2014 the exact tradecraft documented in the Cybereason Cobalt Kitty and ESET Bismuth reports.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.002",
        "T1036.005",
        "T1218",
        "T1027"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.process_hash) as hash from datamodel=Endpoint.Processes where Processes.process_name IN (\"GoogleUpdate.exe\",\"RasTls.exe\",\"rastlsc.exe\",\"mcoemcpy.exe\",\"McVsShld.exe\",\"goopdate.exe\",\"FFViewer.exe\",\"DTHelper.exe\",\"wsc.exe\",\"sysprep.exe\",\"AcroRd32.exe\",\"WerFault.exe\") AND (Processes.process_path=\"*\\\\AppData\\\\*\" OR Processes.process_path=\"*\\\\ProgramData\\\\*\" OR Processes.process_path=\"*\\\\Users\\\\Public\\\\*\" OR Processes.process_path=\"*\\\\Temp\\\\*\" OR Processes.process_path=\"*\\\\PerfLogs\\\\*\" OR Processes.process_path=\"*\\\\Windows\\\\Tasks\\\\*\") AND Processes.user!=\"*$\" by Processes.dest Processes.user Processes.process_name Processes.process_path | `drop_dm_object_name(Processes)` | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "// APT32 DLL side-loading: signed vendor binary executing from user-writable path,\n// then loading a co-resident DLL whose hash differs from the loader.\nlet _apt32_sideload_binaries = dynamic([\n    \"googleupdate.exe\",\"rastls.exe\",\"rastlsc.exe\",\"mcoemcpy.exe\",\n    \"mcvsshld.exe\",\"goopdate.exe\",\"ffviewer.exe\",\"dthelper.exe\",\n    \"wsc.exe\",\"sysprep.exe\",\"acrord32.exe\",\"werfault.exe\"\n]);\nlet _user_writable = dynamic([\n    @\"\\appdata\\\",@\"\\programdata\\\",@\"\\users\\public\\\",\n    @\"\\temp\\\",@\"\\perflogs\\\",@\"\\windows\\tasks\\\"\n]);\nlet LoaderRuns =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where AccountName !endswith \"$\"                          // exclude machine accts\n    | where FileName in~ (_apt32_sideload_binaries)\n    | where tolower(FolderPath) has_any (_user_writable)        // not Program Files\n    | project Timestamp, DeviceId, DeviceName, AccountName,\n              LoaderName = FileName, LoaderPath = FolderPath,\n              LoaderSha256 = SHA256, LoaderCmd = ProcessCommandLine,\n              LoaderProcId = ProcessId;\n// pivot to image loads from the same folder by the same loader process\nDeviceImageLoadEvents\n| where Timestamp > ago(7d)\n| where FileName endswith \".dll\"\n| join kind=inner LoaderRuns on $left.DeviceId == $right.DeviceId,\n                                $left.InitiatingProcessId == $right.LoaderProcId\n| where tolower(FolderPath) == tolower(LoaderPath)              // DLL co-resident\n| where SHA256 != LoaderSha256\n| where not(FolderPath startswith @\"C:\\Windows\\\")\n| project Timestamp, DeviceName, AccountName,\n          LoaderName, LoaderPath, LoaderSha256, LoaderCmd,\n          SideLoadedDll = FileName, SideLoadedDllSha256 = SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// Sentinel-side: Defender doesn't ship DeviceImageLoadEvents to Sentinel by default,\n// so we anchor on the unmistakable behaviour \u2014 a known APT32 carrier binary\n// running from a user-writable directory (4688 / Sysmon EID 1).\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| where (NewProcessName endswith @\"\\GoogleUpdate.exe\"\n      or NewProcessName endswith @\"\\RasTls.exe\"\n      or NewProcessName endswith @\"\\rastlsc.exe\"\n      or NewProcessName endswith @\"\\mcoemcpy.exe\"\n      or NewProcessName endswith @\"\\McVsShld.exe\"\n      or NewProcessName endswith @\"\\goopdate.exe\"\n      or NewProcessName endswith @\"\\FFViewer.exe\"\n      or NewProcessName endswith @\"\\DTHelper.exe\"\n      or NewProcessName endswith @\"\\wsc.exe\"\n      or NewProcessName endswith @\"\\sysprep.exe\"\n      or NewProcessName endswith @\"\\WerFault.exe\")\n| where NewProcessName has_any (@\"\\AppData\\\", @\"\\ProgramData\\\",\n                                 @\"\\Users\\Public\\\", @\"\\Temp\\\",\n                                 @\"\\PerfLogs\\\", @\"\\Windows\\Tasks\\\")\n| project TimeGenerated, Computer, Account,\n          NewProcessName, CommandLine, ParentProcessName,\n          TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "title: APT32 OceanLotus DLL Side-Loading via Signed Binary in User-Writable Path\nid: 7b1a32d0-3c54-4cf2-9f9e-ap32sideload01\nstatus: experimental\ndescription: |\n  Detects APT32 (OceanLotus / Cobalt Kitty / Canvas Cyclone / BISMUTH)\n  DLL side-loading tradecraft where the actor relocates a legitimate signed\n  third-party binary (GoogleUpdate, RasTls, mcoemcpy, goopdate, FFViewer,\n  DTHelper, etc.) to a user-writable directory in order to load a malicious\n  same-folder DLL (KerrDown / Denis / Cobalt Strike stager).\nreferences:\n  - https://attack.mitre.org/groups/G0050/\n  - https://www.cybereason.com/blog/operation-cobalt-kitty-apt\n  - https://www.welivesecurity.com/2020/11/oceanlotus-bismuth-cryptominer/\nauthor: SOC Detection Engineering\ndate: 2026/05/07\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_image:\n    Image|endswith:\n      - '\\GoogleUpdate.exe'\n      - '\\RasTls.exe'\n      - '\\rastlsc.exe'\n      - '\\mcoemcpy.exe'\n      - '\\McVsShld.exe'\n      - '\\goopdate.exe'\n      - '\\FFViewer.exe'\n      - '\\DTHelper.exe'\n      - '\\wsc.exe'\n      - '\\sysprep.exe'\n      - '\\WerFault.exe'\n  selection_path:\n    Image|contains:\n      - '\\AppData\\'\n      - '\\ProgramData\\'\n      - '\\Users\\Public\\'\n      - '\\Temp\\'\n      - '\\PerfLogs\\'\n      - '\\Windows\\Tasks\\'\n  filter_machine_acct:\n    User|endswith: '$'\n  condition: selection_image and selection_path and not filter_machine_acct\nfalsepositives:\n  - Software-installer extraction staging in %TEMP% (rare for these specific binaries)\n  - Out-of-band patcher running from %ProgramData%\nlevel: high\ntags:\n  - attack.execution\n  - attack.persistence\n  - attack.defense_evasion\n  - attack.t1574.002\n  - attack.t1036.005\n  - attack.t1027\n  - attack.g0050",
      "rationale": "Generic side-loading detections fire on every legit installer. This UC narrows to the *specific signed-binary names* documented across APT32 campaigns (Cobalt Kitty: RasTls/mcoemcpy/goopdate/FFViewer/DTHelper; KerrDown: GoogleUpdate; BISMUTH: WerFault/AcroRd32) and additionally requires the binary to be running from a non-Program-Files location \u2014 the structural condition that makes it side-loading",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] APT32 Cobalt Kitty scheduled-task persistence with obfuscated PowerShell payload",
      "description": "APT32's Cobalt Kitty operation registered scheduled tasks under disguised names (e.g. Microsoft\\Windows\\NlaSvc\\WPDBusEnumRoot, Microsoft\\Windows\\Multimedia\\SystemSoundsService) that invoke PowerShell with -EncodedCommand / -nop -w hidden / FromBase64String to launch Cobalt Strike beacons or Denis backdoor stagers. This UC hunts schtasks.exe creating tasks whose action is an obfuscated PowerShell or rundll32 javascript: invocation \u2014 the canonical APT32 persistence shape.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1053.005",
        "T1059.001",
        "T1027",
        "T1036.004",
        "T1218.011"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.parent_process) as parent_cmd from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" AND Processes.process=\"*create*\" AND (Processes.process=\"*EncodedCommand*\" OR Processes.process=\"* -enc *\" OR Processes.process=\"* -ec *\" OR Processes.process=\"*FromBase64String*\" OR Processes.process=\"*-w hidden*\" OR Processes.process=\"*-WindowStyle hidden*\" OR Processes.process=\"*-nop*\" OR Processes.process=\"*IEX*\" OR Processes.process=\"*Invoke-Expression*\" OR Processes.process=\"*javascript:*\" OR Processes.process=\"*rundll32*javascript*\") AND Processes.user!=\"*$\" AND Processes.user!=\"SYSTEM\" by Processes.dest Processes.user Processes.process Processes.parent_process_name | `drop_dm_object_name(Processes)` | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "// APT32 Cobalt-Kitty scheduled-task persistence: schtasks.exe creating a task\n// whose action is an obfuscated PowerShell / rundll32-javascript launcher.\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where AccountName !endswith \"$\"\n| where AccountName !in~ (\"system\",\"local service\",\"network service\")\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has_any (\"/create\",\"-create\",\"/Create\")\n| extend cmd = tolower(ProcessCommandLine)\n| where cmd has_any (\"powershell\",\"pwsh\",\"rundll32\")\n| where cmd has_any (\n    \"-encodedcommand\",\" -enc \",\" -ec \",\"frombase64string\",\n    \"-w hidden\",\"-windowstyle hidden\",\"-nop \",\"-noni\",\n    \"iex \",\"invoke-expression\",\"javascript:\",\n    \"downloadstring\",\"downloadfile\",\"new-object net.webclient\")\n| project Timestamp, DeviceName, AccountName,\n          ProcessCommandLine,\n          ParentBin = InitiatingProcessFileName,\n          ParentCmd = InitiatingProcessCommandLine,\n          ParentParent = InitiatingProcessParentFileName,\n          IsRemote = IsInitiatingProcessRemoteSession\n| order by Timestamp desc",
      "sentinel_kql": "// Sentinel port \u2014 anchored on 4688 process-create.\n// CommandLine collection requires GPO 'Include command line in process creation events'.\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| where AccountName !in~ (\"system\",\"local service\",\"network service\")\n| where NewProcessName endswith @\"\\schtasks.exe\"\n| extend cmd = tolower(CommandLine)\n| where cmd has_any (\"/create\",\"-create\")\n| where cmd has_any (\"powershell\",\"pwsh\",\"rundll32\")\n| where cmd has_any (\n    \"-encodedcommand\",\" -enc \",\" -ec \",\"frombase64string\",\n    \"-w hidden\",\"-windowstyle hidden\",\"-nop \",\"-noni\",\n    \"iex \",\"invoke-expression\",\"javascript:\",\n    \"downloadstring\",\"downloadfile\",\"new-object net.webclient\")\n| project TimeGenerated, Computer, Account, CommandLine,\n          ParentProcessName, TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "title: APT32 OceanLotus Scheduled-Task Persistence with Obfuscated PowerShell Action\nid: 9a4f12bd-2d77-4e9c-abc1-ap32coba1tkit\nstatus: experimental\ndescription: |\n  Detects APT32 (OceanLotus / Cobalt Kitty / BISMUTH) persistence via\n  schtasks.exe registering a task whose action is an obfuscated PowerShell\n  or rundll32-javascript launcher. Documented in the Cybereason Cobalt Kitty\n  report as the actor's primary persistence mechanism for Cobalt Strike\n  beacons and the Denis backdoor.\nreferences:\n  - https://attack.mitre.org/groups/G0050/\n  - https://www.cybereason.com/blog/operation-cobalt-kitty-apt\nauthor: SOC Detection Engineering\ndate: 2026/05/07\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_schtasks:\n    Image|endswith: '\\schtasks.exe'\n  selection_create:\n    CommandLine|contains:\n      - '/create'\n      - '-create'\n  selection_launcher:\n    CommandLine|contains:\n      - 'powershell'\n      - 'pwsh'\n      - 'rundll32'\n  selection_obfuscation:\n    CommandLine|contains:\n      - '-EncodedCommand'\n      - ' -enc '\n      - ' -ec '\n      - 'FromBase64String'\n      - '-w hidden'\n      - '-WindowStyle hidden'\n      - '-nop '\n      - '-noni'\n      - 'IEX '\n      - 'Invoke-Expression'\n      - 'javascript:'\n      - 'DownloadString'\n      - 'New-Object Net.WebClient'\n  filter_machine_acct:\n    User|endswith: '$'\n  condition: selection_schtasks and selection_create and selection_launcher and selection_obfuscation and not filter_machine_acct\nfalsepositives:\n  - Sanctioned admin-driven scheduled jobs that legitimately use -EncodedCommand\n  - Software-deployment systems (SCCM/Intune) launching obfuscated PS via tasks\nlevel: high\ntags:\n  - attack.persistence\n  - attack.execution\n  - attack.defense_evasion\n  - attack.t1053.005\n  - attack.t1059.001\n  - attack.t1027\n  - attack.t1036.004\n  - attack.t1218.011\n  - attack.g0050",
      "rationale": "Rather than alerting on every schtasks /create (way too noisy), this UC requires three concurrent conditions that mirror APT32's documented Cobalt Kitty chain: (a) task creation via schtasks, (b) the action invokes PowerShell or rundll32, (c) the action is obfuscated (EncodedCommand / hidden window / IEX / FromBase64String / javascript:). That tri-condition is what APT32 looks like in their incide",
      "_field_issues": [
        {
          "kind": "unknown_field",
          "field": "IsInitiatingProcessRemoteSession",
          "tables_in_scope": [
            "DeviceProcessEvents"
          ],
          "suggestion": "InitiatingProcessTokenElevation",
          "message": "`IsInitiatingProcessRemoteSession` is not a canonical column of DeviceProcessEvents. Did you mean `InitiatingProcessTokenElevation`?"
        }
      ],
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    }
  ]
}