# intel/.llm_actor_uc_cache/76/76886ba91586871856be02175a48f098eb1c74a4.json
{
  "name": "Qilin",
  "ucs": [
    {
      "title": "[LLM] Qilin (Agenda) GPO-pushed PowerShell Chrome credential harvester",
      "description": "Hunts Qilin's signature mass-credential heist tradecraft (publicly documented by Sophos X-Ops, July 2024): a PowerShell script delivered via Group Policy logon-script / SYSVOL that copies Chrome's `Login Data` SQLite database and `Local State` JSON, extracts the DPAPI-protected `encrypted_key`, and exfils browser passwords from every domain-joined endpoint before encryption. This GPO-as-credential-distribution-channel is highly distinctive and rare outside Qilin/Agenda ransomware intrusions.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1555.003",
        "T1059.001",
        "T1037.001",
        "T1003"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as process values(Processes.process_path) as process_path values(Processes.parent_process) as parent_process values(Processes.parent_process_name) as parent_process_name values(Processes.user) as user dc(Processes.dest) as host_count from datamodel=Endpoint.Processes where (Processes.process_name=\"powershell.exe\" OR Processes.process_name=\"pwsh.exe\") AND (Processes.process=\"*Google\\\\Chrome\\\\User Data\\\\Local State*\" OR Processes.process=\"*Login Data*\" OR (Processes.process=\"*encrypted_key*\" AND Processes.process=\"*FromBase64String*\") OR (Processes.process=\"*ProtectedData*\" AND Processes.process=\"*DPAPI*\") OR Processes.process=\"*AppData\\\\Local\\\\Google\\\\Chrome*Login Data*\") AND (Processes.parent_process_name=\"gpscript.exe\" OR Processes.parent_process_name=\"svchost.exe\" OR Processes.parent_process_name=\"wmiprvse.exe\" OR Processes.parent_process_name=\"taskeng.exe\" OR Processes.parent_process=\"*\\\\SYSVOL\\\\*\" OR Processes.parent_process_name=\"schtasks.exe\") by Processes.dest Processes.user Processes.process_name Processes.process Processes.parent_process_name Processes.parent_process | `drop_dm_object_name(Processes)` | where host_count >= 2 | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "let LookbackDays = 7d;\nlet ChromeArtifacts = dynamic([@\"\\\\Google\\\\Chrome\\\\User Data\\\\Local State\", @\"\\\\Google\\\\Chrome\\\\User Data\\\\Default\\\\Login Data\", @\"\\\\AppData\\\\Local\\\\Google\\\\Chrome\", \"Login Data\", \"Local State\"]);\nlet DpapiTokens   = dynamic([\"ProtectedData\", \"Unprotect\", \"DPAPI\", \"encrypted_key\", \"AesGcm\", \"System.Security.Cryptography\"]);\nlet GpoParents    = dynamic([\"gpscript.exe\", \"svchost.exe\", \"wmiprvse.exe\", \"taskeng.exe\", \"schtasks.exe\", \"services.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(LookbackDays)\n| where FileName in~ (\"powershell.exe\", \"pwsh.exe\")\n| where AccountName !endswith \"$\" or InitiatingProcessFolderPath has @\"\\SYSVOL\\\"   // keep SYSTEM context if it came from a GPO logon script\n| where ProcessCommandLine has_any (ChromeArtifacts)\n      or (ProcessCommandLine has_any (DpapiTokens) and ProcessCommandLine has \"Chrome\")\n| where InitiatingProcessFileName in~ (GpoParents)\n      or InitiatingProcessCommandLine has_any (@\"\\SYSVOL\\\", @\"\\NETLOGON\\\", \"gpt.ini\", \"GroupPolicy\")\n      or InitiatingProcessFolderPath has @\"\\SYSVOL\\\"\n| project Timestamp, DeviceName, AccountName,\n          ParentImage = InitiatingProcessFolderPath,\n          ParentCmd   = InitiatingProcessCommandLine,\n          ChildImage  = FolderPath,\n          ChildCmd    = ProcessCommandLine,\n          IsRemoteSession = InitiatingProcessTokenElevation,\n          SHA256\n| summarize FirstSeen = min(Timestamp), LastSeen = max(Timestamp),\n            HostCount = dcount(DeviceName),\n            Devices   = make_set(DeviceName, 50),\n            SampleCmd = any(ChildCmd)\n            by ParentImage, ChildImage\n| where HostCount >= 2     // GPO-distributed payload \u2014 fans out across the estate\n| order by FirstSeen desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet ChromeArtifacts = dynamic([@\"\\\\Google\\\\Chrome\\\\User Data\\\\Local State\", @\"\\\\Google\\\\Chrome\\\\User Data\\\\Default\\\\Login Data\", \"Login Data\", \"Local State\", @\"\\\\AppData\\\\Local\\\\Google\\\\Chrome\"]);\nlet DpapiTokens     = dynamic([\"ProtectedData\", \"Unprotect\", \"DPAPI\", \"encrypted_key\", \"AesGcm\"]);\nlet GpoParents      = dynamic([\"gpscript.exe\", \"svchost.exe\", \"wmiprvse.exe\", \"taskeng.exe\", \"schtasks.exe\", \"services.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventID == 4688                              // process create \u2014 requires \"Include command line\" GPO\n| where NewProcessName endswith \"\\\\powershell.exe\" or NewProcessName endswith \"\\\\pwsh.exe\"\n| where CommandLine has_any (ChromeArtifacts)\n      or (CommandLine has_any (DpapiTokens) and CommandLine has \"Chrome\")\n| where ParentProcessName has_any (GpoParents)\n      or CommandLine has_any (@\"\\SYSVOL\\\", @\"\\NETLOGON\\\")\n| project TimeGenerated, Computer, Account, NewProcessName, CommandLine, ParentProcessName, TokenElevationType\n| summarize FirstSeen = min(TimeGenerated), LastSeen = max(TimeGenerated),\n            HostCount = dcount(Computer),\n            Hosts     = make_set(Computer, 50),\n            SampleCmd = any(CommandLine)\n            by ParentProcessName, NewProcessName\n| where HostCount >= 2\n| order by FirstSeen desc",
      "sigma_yaml": "title: Qilin (Agenda) GPO-Distributed Chrome Credential Harvester\nid: 5b3d6e1c-7c2b-4d59-aef0-5e7d1f6a91d4\nstatus: experimental\ndescription: Detects PowerShell launched from a GPO/logon-script context that reads Chrome's Login Data or Local State, characteristic of Qilin/Agenda mass-credential staging.\nauthor: SOC Detection Engineering\ndate: 2026/05/06\nreferences:\n  - https://news.sophos.com/en-us/2024/08/22/qilin-ransomware-caught-stealing-credentials-stored-in-google-chrome/\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_proc:\n    Image|endswith:\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n  selection_chrome:\n    CommandLine|contains:\n      - '\\Google\\Chrome\\User Data\\Local State'\n      - '\\Google\\Chrome\\User Data\\Default\\Login Data'\n      - 'Login Data'\n      - 'Local State'\n      - 'encrypted_key'\n  selection_parent_or_path:\n    - ParentImage|endswith:\n        - '\\gpscript.exe'\n        - '\\wmiprvse.exe'\n        - '\\taskeng.exe'\n        - '\\schtasks.exe'\n    - CommandLine|contains:\n        - '\\SYSVOL\\'\n        - '\\NETLOGON\\'\n    - ParentCommandLine|contains:\n        - '\\SYSVOL\\'\n        - 'GroupPolicy'\n  condition: selection_proc and selection_chrome and selection_parent_or_path\nfields:\n  - Image\n  - CommandLine\n  - ParentImage\n  - ParentCommandLine\n  - User\n  - Computer\nfalsepositives:\n  - Authorised internal credential-audit / browser-migration tooling distributed via GPO\nlevel: high\ntags:\n  - attack.credential_access\n  - attack.t1555.003\n  - attack.execution\n  - attack.t1059.001\n  - attack.persistence\n  - attack.t1037.001",
      "rationale": "Generic 'PowerShell touches Chrome' detections are noisy. What makes this Qilin-specific is the *delivery vector*: the credential-stealer is pushed estate-wide via a Group Policy logon script (parent = gpscript.exe / SYSVOL path / wmiprvse.exe under SYSTEM) and fans out to many endpoints simultaneously. This combination \u2014 DPAPI/Chrome-credential code paths in PowerShell, GPO/SYSVOL parent context,"
    },
    {
      "title": "[LLM] Qilin (Agenda) Linux/ESXi Rust encryptor pre-encryption VM-kill chain",
      "description": "Hunts Qilin's ESXi/Linux variant \u2014 a Rust-compiled binary that enumerates virtual machines via `esxcli`/`vim-cmd`, force-kills running guests, then walks `/vmfs/volumes/` encrypting `.vmdk`/`.vmx`/`.vswp` files. The combination of `esxcli vm process kill` followed within minutes by mass file rename/encryption under `/vmfs/volumes` from an unsigned binary in `/tmp` is a Qilin/Agenda hallmark and not a normal ESXi admin workflow.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1486",
        "T1489",
        "T1059.004",
        "T1490",
        "T1078"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as process values(Processes.user) as user values(Processes.parent_process_name) as parent from datamodel=Endpoint.Processes where Processes.os=\"Linux\" AND ((Processes.process_name=\"esxcli\" AND Processes.process=\"*vm process kill*\" AND Processes.process=\"*--type=force*\") OR (Processes.process_name=\"vim-cmd\" AND Processes.process=\"*vmsvc/getallvms*\") OR (Processes.process=\"*/tmp/*\" AND (Processes.process=\"*-yes*\" OR Processes.process=\"*--debug*\" OR Processes.process=\"*--logs*\" OR Processes.process=\"*-path*/vmfs/volumes*\")) OR (Processes.process_name=\"chmod\" AND Processes.process=\"*+x*\" AND Processes.process=\"*/tmp/*\")) by Processes.dest Processes.process_name Processes.process Processes.parent_process_name Processes.user _time | `drop_dm_object_name(Processes)` | eventstats dc(process_name) as distinct_actions by dest | where distinct_actions >= 2 | sort 0 -firstTime | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "let LookbackHours = 6h;\nlet KillVm   = DeviceProcessEvents\n    | where Timestamp > ago(LookbackHours)\n    | where DeviceName has_any (\"esxi\", \"vmhost\") or InitiatingProcessFolderPath has \"/bin/\" or FolderPath has \"/bin/\"\n    | where (FileName =~ \"esxcli\" and ProcessCommandLine has_all (\"vm\", \"process\", \"kill\"))\n         or (FileName =~ \"vim-cmd\" and ProcessCommandLine has \"vmsvc/getallvms\");\nlet Encryptor = DeviceProcessEvents\n    | where Timestamp > ago(LookbackHours)\n    | where (FolderPath startswith \"/tmp/\" or FolderPath startswith \"/var/tmp/\" or FolderPath startswith \"/vmfs/volumes/\")\n    | where ProcessCommandLine has_any (\"-yes\", \"--debug\", \"--logs\", \"--no-vm-kill\", \"-path\", \"--password\")\n         or ProcessCommandLine has \"/vmfs/volumes\";\nKillVm\n| project KillTime = Timestamp, DeviceId, DeviceName, KillCmd = ProcessCommandLine\n| join kind=inner (\n    Encryptor\n    | project EncTime = Timestamp, DeviceId, EncImage = FolderPath, EncCmd = ProcessCommandLine, EncSHA256 = SHA256\n  ) on DeviceId\n| where EncTime between (KillTime .. KillTime + 30m)\n| project DeviceName, KillTime, KillCmd, EncTime, EncImage, EncCmd, EncSHA256,\n          DelayMin = datetime_diff('minute', EncTime, KillTime)\n| order by EncTime desc",
      "sentinel_kql": "let LookbackHours = 6h;\nlet KillVm = ImProcessCreate\n    | where TimeGenerated > ago(LookbackHours)\n    | where DvcOs has_any (\"linux\", \"esxi\") or DvcHostname has_any (\"esxi\", \"vmhost\")\n    | where (TargetProcessName endswith \"esxcli\" and TargetProcessCommandLine has_all (\"vm\", \"process\", \"kill\"))\n         or (TargetProcessName endswith \"vim-cmd\" and TargetProcessCommandLine has \"vmsvc/getallvms\")\n    | project KillTime = TimeGenerated, DvcHostname, KillCmd = TargetProcessCommandLine, ActorUsername;\nlet Encryptor = ImProcessCreate\n    | where TimeGenerated > ago(LookbackHours)\n    | where TargetProcessName matches regex @\"^/(tmp|var/tmp|vmfs/volumes)/.+\"\n    | where TargetProcessCommandLine has_any (\"-yes\", \"--debug\", \"--logs\", \"--no-vm-kill\", \"-path\", \"--password\")\n         or TargetProcessCommandLine has \"/vmfs/volumes\"\n    | project EncTime = TimeGenerated, DvcHostname, EncImage = TargetProcessName, EncCmd = TargetProcessCommandLine, EncSHA256 = SHA256;\nKillVm\n| join kind=inner Encryptor on DvcHostname\n| where EncTime between (KillTime .. KillTime + 30m)\n| project DvcHostname, ActorUsername, KillTime, KillCmd, EncTime, EncImage, EncCmd, EncSHA256,\n          DelayMin = datetime_diff('minute', EncTime, KillTime)\n| order by EncTime desc",
      "sigma_yaml": "",
      "rationale": "Qilin (alongside affiliates of Akira and BlackBasta) is one of the few crews that ships a *Rust-compiled* ESXi encryptor with the documented `-yes / --logs / --debug / --no-vm-kill / -path` flag set and a built-in VM-killer that calls `esxcli vm process kill --type=force` before encrypting `/vmfs/volumes/`. The two-stage chain (esxcli VM kill \u2192 unsigned `/tmp` binary writing under `/vmfs/volumes` "
    }
  ]
}