# intel/.llm_actor_uc_cache/ee/ee60ea9bfaa55839b81a98abc2133d23af7f42d8.json
{
  "name": "Stealth Falcon",
  "ucs": [
    {
      "title": "[LLM] Stealth Falcon BITS-job C2 / exfiltration via PowerShell or bitsadmin.exe",
      "description": "Hunts Stealth Falcon's signature backdoor (documented by ESET as Win32/StealthFalcon) which abuses the Background Intelligent Transfer Service to stage covert HTTPS C2 and exfiltrate stolen data \u2014 either via direct bitsadmin.exe /transfer invocations or via PowerShell loading BITS COM classes (Start-BitsTransfer / BackgroundCopyManager). This is the actor's tradecraft fingerprint, not generic LOLBin abuse.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "high",
      "techs": [
        "T1071.001",
        "T1059.001",
        "T1041",
        "T1573.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where (Processes.process_name=\"bitsadmin.exe\" AND (Processes.process=\"*/transfer*\" OR Processes.process=\"*/addfile*\" OR Processes.process=\"*/setnotifycmdline*\" OR Processes.process=\"*/setminretrydelay*\" OR Processes.process=\"*/setcustomheaders*\" OR Processes.process=\"*http*\")) OR (Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\") AND (Processes.process=\"*Start-BitsTransfer*\" OR Processes.process=\"*BackgroundCopyManager*\" OR Processes.process=\"*BITSAdmin*\" OR Processes.process=\"*Microsoft.BackgroundIntelligentTransfer*\")) by host Processes.process_name Processes.user Processes.parent_process_name _time span=1m | where user!=\"SYSTEM\" AND parent!=\"WindowsUpdateBox.exe\" AND parent!=\"TiWorker.exe\" AND parent!=\"MoUsoCoreWorker.exe\" | `drop_dm_object_name(Processes)` | join type=inner host [| tstats summariesonly=t values(All_Traffic.dest) as dest values(All_Traffic.dest_port) as dest_port from datamodel=Network_Traffic.All_Traffic where All_Traffic.app=\"bitsadmin.exe\" OR All_Traffic.app=\"svchost.exe\" by host _time span=5m] | where firstTime <= _time AND _time <= lastTime+300 | sort -firstTime",
      "kql": "// Stealth Falcon \u2014 BITS-based C2 (ESET Win32/StealthFalcon, 2019)\nlet lookback = 7d;\nlet ms_signed_parents = dynamic([\"WindowsUpdateBox.exe\",\"TiWorker.exe\",\"MoUsoCoreWorker.exe\",\"sdiagnhost.exe\",\"OneDriveSetup.exe\"]);\nlet bits_proc = DeviceProcessEvents\n    | where Timestamp > ago(lookback)\n    | where AccountName !endswith \"$\"   // exclude machine accounts\n    | where (FileName =~ \"bitsadmin.exe\"\n             and ProcessCommandLine has_any (\"/transfer\",\"/addfile\",\"/setnotifycmdline\",\"/setcustomheaders\",\"/setminretrydelay\"))\n         or (FileName in~ (\"powershell.exe\",\"pwsh.exe\")\n             and ProcessCommandLine has_any (\"Start-BitsTransfer\",\"BackgroundCopyManager\",\"Microsoft.BackgroundIntelligentTransfer\",\"BITSAdmin\"))\n    | where InitiatingProcessFileName !in~ (ms_signed_parents)\n    | project ProcTime = Timestamp, DeviceId, DeviceName, AccountName,\n              FileName, ProcessCommandLine, InitiatingProcessFileName,\n              InitiatingProcessCommandLine;\n// Correlate to a same-host outbound HTTPS connection from BITS host process within 10 min\nbits_proc\n| join kind=inner (\n    DeviceNetworkEvents\n    | where Timestamp > ago(lookback)\n    | where RemoteIPType == \"Public\"\n    | where RemotePort in (80, 443, 8080, 8443)\n    | where InitiatingProcessFileName in~ (\"svchost.exe\",\"bitsadmin.exe\",\"powershell.exe\",\"pwsh.exe\")\n    | project NetTime = Timestamp, DeviceId, RemoteIP, RemoteUrl, RemotePort,\n              NetProc = InitiatingProcessFileName, NetCmd = InitiatingProcessCommandLine\n  ) on DeviceId\n| where NetTime between (ProcTime .. ProcTime + 10m)\n| summarize Connections = count(),\n            DistinctRemotes = dcount(RemoteIP),\n            SampleUrls = make_set(RemoteUrl, 8),\n            SampleIPs  = make_set(RemoteIP, 8)\n            by DeviceName, AccountName, FileName, ProcessCommandLine, InitiatingProcessFileName, ProcTime\n| where DistinctRemotes >= 1\n| order by ProcTime desc",
      "sentinel_kql": "// Stealth Falcon \u2014 BITS-based C2 (Win32/StealthFalcon)\nlet lookback = 7d;\nlet ms_signed_parents = dynamic([\"WindowsUpdateBox.exe\",\"TiWorker.exe\",\"MoUsoCoreWorker.exe\",\"sdiagnhost.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(lookback)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| where (NewProcessName endswith \"\\\\bitsadmin.exe\"\n         and CommandLine has_any (\"/transfer\",\"/addfile\",\"/setnotifycmdline\",\"/setcustomheaders\",\"/setminretrydelay\"))\n     or ((NewProcessName endswith \"\\\\powershell.exe\" or NewProcessName endswith \"\\\\pwsh.exe\")\n         and CommandLine has_any (\"Start-BitsTransfer\",\"BackgroundCopyManager\",\"Microsoft.BackgroundIntelligentTransfer\",\"BITSAdmin\"))\n| extend ParentLeaf = tostring(split(ParentProcessName, \"\\\\\")[-1])\n| where ParentLeaf !in~ (ms_signed_parents)\n| project TimeGenerated, Computer, Account, NewProcessName, CommandLine,\n          ParentProcessName, TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Stealth Falcon BITS Job Abuse for Covert C2\nid: 4f7c2a18-8b65-4c3a-9d11-bits-stealthfalcon\nstatus: experimental\ndescription: Detects Stealth Falcon backdoor C2/exfil tradecraft via bitsadmin.exe /transfer or PowerShell BITS COM class invocation\nreferences:\n  - https://www.welivesecurity.com/2019/09/09/backdoor-stealth-falcon-group/\n  - https://attack.mitre.org/groups/G0038/\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.command_and_control\n  - attack.t1071.001\n  - attack.t1059.001\n  - attack.g0038\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_bitsadmin:\n    Image|endswith: '\\bitsadmin.exe'\n    CommandLine|contains:\n      - '/transfer'\n      - '/addfile'\n      - '/setnotifycmdline'\n      - '/setcustomheaders'\n      - '/setminretrydelay'\n  selection_pwsh_bits:\n    Image|endswith:\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n    CommandLine|contains:\n      - 'Start-BitsTransfer'\n      - 'BackgroundCopyManager'\n      - 'Microsoft.BackgroundIntelligentTransfer'\n  filter_legit:\n    ParentImage|endswith:\n      - '\\WindowsUpdateBox.exe'\n      - '\\TiWorker.exe'\n      - '\\MoUsoCoreWorker.exe'\n      - '\\sdiagnhost.exe'\n  condition: (selection_bitsadmin or selection_pwsh_bits) and not filter_legit\nfalsepositives:\n  - Software-distribution agents using BITS for content delivery\n  - Legitimate Windows Update / SCCM workflows\nlevel: high",
      "rationale": "Stealth Falcon's backdoor (G0038, ESET 2019 'Win32/StealthFalcon') is unusual among APTs in that it deliberately uses the BITS service for C2 staging and exfiltration \u2014 most APT C2 uses raw sockets or wininet. Pairing the BITS-job artefact (bitsadmin.exe /transfer flags or PowerShell BITS COM classes) with an outbound HTTPS connection from the same host within 10 minutes captures this exact tradec",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] Stealth Falcon scheduled-task PowerShell collecting Windows Vault & browser credential stores",
      "description": "Hunts Stealth Falcon's post-compromise credential harvesting pattern: a scheduled task (T1053.005) launches PowerShell (T1059.001) which then reads the Windows Credential Manager Vault directory or Chromium 'Login Data' SQLite database (T1555.003/T1555.004) and stages output for exfil. This combined sequence is characteristic of the actor's spyware against UAE journalists/dissidents.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1053.005",
        "T1059.001",
        "T1555.003",
        "T1555.004",
        "T1041"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as ps_cmd values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\") AND Processes.parent_process_name IN (\"svchost.exe\",\"taskeng.exe\",\"schtasks.exe\",\"WmiPrvSE.exe\") by host Processes.process_id _time span=1m | `drop_dm_object_name(Processes)` | join type=inner host process_id [| tstats summariesonly=t values(Filesystem.file_path) as touched_files from datamodel=Endpoint.Filesystem where (Filesystem.file_path=\"*\\\\Microsoft\\\\Vault\\\\*\" OR Filesystem.file_path=\"*\\\\Microsoft\\\\Credentials\\\\*\" OR Filesystem.file_name=\"Login Data\" OR Filesystem.file_name=\"Cookies\" OR Filesystem.file_name=\"Web Data\" OR Filesystem.file_path=\"*\\\\Mozilla\\\\Firefox\\\\Profiles\\\\*key*.db*\" OR Filesystem.file_path=\"*\\\\Mozilla\\\\Firefox\\\\Profiles\\\\*logins.json*\") by host Filesystem.process_id] | where user!=\"SYSTEM\" | sort -firstTime",
      "kql": "// Stealth Falcon \u2014 scheduled-task PowerShell harvesting Vault + browser creds\nlet lookback = 7d;\nlet sched_parents = dynamic([\"svchost.exe\",\"taskeng.exe\",\"schtasks.exe\",\"WmiPrvSE.exe\"]);\nlet cred_paths = dynamic([\n    @\"\\Microsoft\\Vault\\\",\n    @\"\\Microsoft\\Credentials\\\",\n    @\"\\Microsoft\\Protect\\\",\n    @\"\\Google\\Chrome\\User Data\\\",\n    @\"\\Microsoft\\Edge\\User Data\\\",\n    @\"\\Mozilla\\Firefox\\Profiles\\\",\n    @\"\\BraveSoftware\\Brave-Browser\\User Data\\\"\n]);\nlet cred_files = dynamic([\"Login Data\",\"Cookies\",\"Web Data\",\"logins.json\",\"key4.db\",\"key3.db\",\"Local State\"]);\n// 1. PowerShell launched by a scheduled-task / WMI host\nlet ps_via_sched = DeviceProcessEvents\n    | where Timestamp > ago(lookback)\n    | where AccountName !endswith \"$\"\n    | where FileName in~ (\"powershell.exe\",\"pwsh.exe\")\n    | where InitiatingProcessFileName in~ (sched_parents)\n         or (InitiatingProcessFileName =~ \"svchost.exe\"\n             and InitiatingProcessCommandLine has_any (\"Schedule\",\"netsvcs\"))\n    | project ProcTime = Timestamp, DeviceId, DeviceName, AccountName,\n              PsId = ProcessId, PsCmd = ProcessCommandLine,\n              Parent = InitiatingProcessFileName,\n              ParentCmd = InitiatingProcessCommandLine;\n// 2. Same PowerShell PID then touches a credential store\nps_via_sched\n| join kind=inner (\n    DeviceFileEvents\n    | where Timestamp > ago(lookback)\n    | where InitiatingProcessFileName in~ (\"powershell.exe\",\"pwsh.exe\")\n    | where FolderPath has_any (cred_paths) or FileName in~ (cred_files)\n    | project FileTime = Timestamp, DeviceId, InitiatingProcessId,\n              FolderPath, FileName, ActionType\n  ) on DeviceId, $left.PsId == $right.InitiatingProcessId\n| where FileTime between (ProcTime .. ProcTime + 30m)\n| summarize CredArtefacts = make_set(strcat(FolderPath, FileName), 25),\n            ArtefactCount = dcount(strcat(FolderPath, FileName)),\n            FirstAccess = min(FileTime), LastAccess = max(FileTime)\n            by DeviceName, AccountName, Parent, ParentCmd, PsCmd, ProcTime\n| where ArtefactCount >= 2   // single-file touches are noisy; chained access is the signal\n| order by ProcTime desc",
      "sentinel_kql": "// Stealth Falcon \u2014 scheduled-task PowerShell + credential-store access\nlet lookback = 7d;\nlet sched_parents = dynamic([\"\\\\svchost.exe\",\"\\\\taskeng.exe\",\"\\\\schtasks.exe\",\"\\\\WmiPrvSE.exe\"]);\nlet cred_paths = dynamic([\n    @\"\\Microsoft\\Vault\\\",\n    @\"\\Microsoft\\Credentials\\\",\n    @\"\\Google\\Chrome\\User Data\\\",\n    @\"\\Microsoft\\Edge\\User Data\\\",\n    @\"\\Mozilla\\Firefox\\Profiles\\\"\n]);\nlet cred_files = dynamic([\"Login Data\",\"Cookies\",\"Web Data\",\"logins.json\",\"key4.db\",\"Local State\"]);\nlet ps_via_sched = SecurityEvent\n    | where TimeGenerated > ago(lookback)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | where NewProcessName endswith \"\\\\powershell.exe\" or NewProcessName endswith \"\\\\pwsh.exe\"\n    | where ParentProcessName has_any (sched_parents)\n    | project ProcTime = TimeGenerated, Computer, Account,\n              PsPid = NewProcessId, PsCmd = CommandLine,\n              Parent = ParentProcessName;\nps_via_sched\n| join kind=inner (\n    ImFileEvent\n    | where TimeGenerated > ago(lookback)\n    | where ActingProcessName endswith \"powershell.exe\" or ActingProcessName endswith \"pwsh.exe\"\n    | where TargetFilePath has_any (cred_paths) or TargetFileName in~ (cred_files)\n    | project FileTime = TimeGenerated, Computer = DvcHostname,\n              ActingProcessId, TargetFilePath, TargetFileName, EventType\n  ) on Computer\n| where FileTime between (ProcTime .. ProcTime + 30m)\n| summarize CredArtefacts = make_set(strcat(TargetFilePath, TargetFileName), 25),\n            ArtefactCount = dcount(strcat(TargetFilePath, TargetFileName)),\n            FirstAccess = min(FileTime), LastAccess = max(FileTime)\n            by Computer, Account, Parent, PsCmd, ProcTime\n| where ArtefactCount >= 2\n| order by ProcTime desc",
      "sigma_yaml": "",
      "rationale": "Stealth Falcon's targeted spyware against Emirati journalists/activists relies on persistent scheduled tasks (T1053.005) that drop a PowerShell-based collector (T1059.001) which then loots both the Windows Credential Vault (T1555.004) and browser-saved credentials (T1555.003). Catching scheduled-task-spawned PowerShell in isolation is noisy; chaining it via PID/host to file access against \u22652 disti"
    }
  ]
}