# intel/.llm_actor_uc_cache/58/582b408affa3060793a0d3397ad6edec603f1eda.json
{
  "name": "APT34",
  "ucs": [
    {
      "title": "[LLM] APT34 (OilRig) DNS-tunneling C2 \u2014 high-cardinality subdomains under attacker-owned parent (DNSpionage / ALMA Communicator / BONDUPDATER / He",
      "description": "Hunts the signature OilRig C2 channel: implants (DNSpionage, ALMA Communicator, Helminth, BONDUPDATER, POWRUNER) encode beacon traffic into long, randomised subdomain labels and exfiltrate via DNS A/TXT lookups. We aggregate per-host DNS responses by registrable parent and flag bursts of many unique long subdomains \u2014 a behaviour OilRig has used since at least 2016 and that distinguishes them from generic PowerShell/Cobalt-Strike crews.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "high",
      "techs": [
        "T1071.004",
        "T1048.003",
        "T1008",
        "T1132"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count values(DNS.query) as queries dc(DNS.query) as unique_queries values(DNS.src) as src values(DNS.record_type) as record_types from datamodel=Network_Resolution where DNS.message_type=\"QUERY\" NOT (DNS.query IN (\"*.microsoft.com\",\"*.windows.com\",\"*.akamai.net\",\"*.azureedge.net\",\"*.cloudfront.net\",\"*.googleapis.com\",\"*.office.com\",\"*.office365.com\",\"*.windowsupdate.com\",\"*.symantec.com\",\"*.skype.com\")) by DNS.src host _time span=1h DNS.query | rex field=DNS.query \"(?<subdomain>.*?)\\.(?<parent_domain>[^.]+\\.[^.]+)$\" | eval sub_len=len(subdomain) | stats count as total_queries dc(subdomain) as unique_subs avg(sub_len) as avg_sub_len max(sub_len) as max_sub_len values(record_types) as record_types values(DNS.query) as sample_queries by host DNS.src parent_domain _time | where unique_subs>=30 AND avg_sub_len>=25 AND (match(record_types,\"TXT\") OR max_sub_len>=50) | sort - unique_subs",
      "kql": "// APT34 / OilRig DNS-tunnel beacon \u2014 long randomised subdomain bursts under a single registrable parent\nlet _allowed_parents = dynamic([\"microsoft.com\",\"windows.com\",\"akamaiedge.net\",\"akamai.net\",\"azureedge.net\",\"office.com\",\"office365.com\",\"trafficmanager.net\",\"googleapis.com\",\"cloudfront.net\",\"windowsupdate.com\",\"skype.com\",\"live.com\",\"msftncsi.com\"]);\nDeviceEvents\n| where Timestamp > ago(7d)\n| where ActionType == \"DnsQueryResponse\"\n| extend QueryName = tolower(tostring(parse_json(AdditionalFields).QueryName))\n| where isnotempty(QueryName) and QueryName has \".\"\n| extend Labels = split(QueryName, \".\")\n| where array_length(Labels) >= 3                              // need at least sub.parent.tld\n| extend ParentDomain = strcat(tostring(Labels[array_length(Labels)-2]), \".\", tostring(Labels[array_length(Labels)-1]))\n| where ParentDomain !in (_allowed_parents)\n| extend Subdomain = substring(QueryName, 0, strlen(QueryName) - strlen(ParentDomain) - 1)\n| where strlen(Subdomain) >= 20                                // OilRig subdomains carry encoded payloads\n| summarize UniqueSubdomains = dcount(Subdomain),\n            QueryCount      = count(),\n            AvgSubLen       = avg(strlen(Subdomain)),\n            MaxSubLen       = max(strlen(Subdomain)),\n            SampleQueries   = make_set(QueryName, 5),\n            Processes       = make_set(InitiatingProcessFileName, 10),\n            FirstSeen       = min(Timestamp),\n            LastSeen        = max(Timestamp)\n          by DeviceName, ParentDomain, bin(Timestamp, 1h)\n| where UniqueSubdomains >= 30 and AvgSubLen >= 25             // OilRig empirical baseline \u2014 random labels, not CDN hashes\n| where Processes !has \"chrome.exe\" and Processes !has \"msedge.exe\" and Processes !has \"firefox.exe\"\n| order by UniqueSubdomains desc",
      "sentinel_kql": "// APT34 / OilRig DNS-tunnel beacon \u2014 Sentinel via Windows DNS server logs + ASIM\nlet _allowed_parents = dynamic([\"microsoft.com\",\"windows.com\",\"akamaiedge.net\",\"akamai.net\",\"azureedge.net\",\"office.com\",\"office365.com\",\"trafficmanager.net\",\"googleapis.com\",\"cloudfront.net\",\"windowsupdate.com\",\"skype.com\",\"live.com\",\"msftncsi.com\"]);\nlet _dns = union isfuzzy=true\n    ( ImDnsActivity\n        | where TimeGenerated > ago(7d)\n        | project TimeGenerated, SrcIpAddr, SrcHostname, DnsQuery = tolower(DnsQuery), DnsQueryTypeName ),\n    ( DnsEvents\n        | where TimeGenerated > ago(7d)\n        | project TimeGenerated, SrcIpAddr = ClientIP, SrcHostname = Computer, DnsQuery = tolower(Name), DnsQueryTypeName = QueryTypeName );\n_dns\n| where isnotempty(DnsQuery) and DnsQuery has \".\"\n| extend Labels = split(DnsQuery, \".\")\n| where array_length(Labels) >= 3\n| extend ParentDomain = strcat(tostring(Labels[array_length(Labels)-2]), \".\", tostring(Labels[array_length(Labels)-1]))\n| where ParentDomain !in (_allowed_parents)\n| extend Subdomain = substring(DnsQuery, 0, strlen(DnsQuery) - strlen(ParentDomain) - 1)\n| where strlen(Subdomain) >= 20\n| summarize UniqueSubdomains = dcount(Subdomain),\n            QueryCount      = count(),\n            AvgSubLen       = avg(strlen(Subdomain)),\n            MaxSubLen       = max(strlen(Subdomain)),\n            QueryTypes      = make_set(DnsQueryTypeName),\n            SampleQueries   = make_set(DnsQuery, 5),\n            FirstSeen       = min(TimeGenerated),\n            LastSeen        = max(TimeGenerated)\n          by SrcHostname, SrcIpAddr, ParentDomain, bin(TimeGenerated, 1h)\n| where UniqueSubdomains >= 30 and AvgSubLen >= 25\n| order by UniqueSubdomains desc",
      "sigma_yaml": "",
      "rationale": "OilRig is one of the very few APTs whose entire toolset family (DNSpionage, ALMA, Helminth, BONDUPDATER, POWRUNER) is built around DNS-tunnel C2 \u2014 Cisco Talos, Palo Alto Unit 42 and FireEye/Mandiant have all documented the same pattern of randomised \u226520-char subdomain labels under attacker parents (e.g. *.poison-frog[.]club, *.windowsdefenderr[.]com, *.withyourface[.]com). Generic Cobalt-Strike ac"
    },
    {
      "title": "[LLM] APT34 (OilRig / Cobalt Gypsy) IIS web-shell \u2014 w3wp.exe spawning interactive shells / discovery (TwoFace / RGDoor / Karkoff)",
      "description": "Catches OilRig's signature post-exploitation pattern on Internet-facing IIS / Exchange: their TwoFace, RGDoor and Karkoff web shells run inside w3wp.exe and routinely invoke cmd.exe, powershell.exe, whoami, net.exe, systeminfo and quser for host triage and dropper staging. Any w3wp child shell is suspicious; APT34's workflow chains discovery + scheduled-task persistence within seconds of first shell.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1505.003",
        "T1059.001",
        "T1059.003",
        "T1087.001",
        "T1087.002",
        "T1082",
        "T1053.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdlines values(Processes.process_name) as child_names values(Processes.user) as users values(Processes.process_id) as pids from datamodel=Endpoint.Processes where Processes.parent_process_name=\"w3wp.exe\" Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"powershell_ise.exe\",\"pwsh.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"systeminfo.exe\",\"quser.exe\",\"hostname.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"schtasks.exe\",\"reg.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"nltest.exe\") by host Processes.user Processes.parent_process_id _time span=10m | rename Processes.* as * | where NOT match(user,\"(?i)^(IIS APPPOOL|NT AUTHORITY|IUSR_)\") OR match(cmdlines,\"(?i)(whoami|systeminfo|net\\s+(user|group|localgroup)|quser|nltest|schtasks\\s+/create)\") | sort - lastTime",
      "kql": "// APT34 / OilRig \u2014 IIS worker process spawning interactive shells or recon binaries (TwoFace / RGDoor / Karkoff family)\nlet _recon_or_shell = dynamic([\"cmd.exe\",\"powershell.exe\",\"powershell_ise.exe\",\"pwsh.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"systeminfo.exe\",\"quser.exe\",\"hostname.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"schtasks.exe\",\"reg.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"nltest.exe\",\"wmic.exe\",\"netstat.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName =~ \"w3wp.exe\"                       // IIS / Exchange OWA / ECP worker\n| where FileName in~ (_recon_or_shell)\n| extend CmdLower = tolower(ProcessCommandLine)\n| extend OilRigRecon = CmdLower has_any (\"whoami\",\"systeminfo\",\"net user\",\"net group\",\"net localgroup\",\"net view\",\"nltest /domain_trusts\",\"quser\",\"tasklist /v\",\"schtasks /create\",\"-encodedcommand\",\"-enc \",\"downloadstring\",\"iex\",\"invoke-expression\",\"invoke-webrequest\",\"certutil -urlcache\",\"bitsadmin /transfer\")\n| project Timestamp, DeviceName,\n          IISAppPoolUser   = InitiatingProcessAccountName,\n          ParentCmd        = InitiatingProcessCommandLine,\n          Child            = FileName,\n          ChildCmd         = ProcessCommandLine,\n          OilRigRecon,\n          ChildSHA256      = SHA256,\n          IsRemoteSession  = InitiatingProcessAccountName\n| order by Timestamp desc",
      "sentinel_kql": "// APT34 / OilRig \u2014 IIS w3wp.exe spawning shells / recon (Sentinel via 4688)\nlet _recon_or_shell = dynamic([\"cmd.exe\",\"powershell.exe\",\"powershell_ise.exe\",\"pwsh.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"systeminfo.exe\",\"quser.exe\",\"hostname.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"schtasks.exe\",\"reg.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"nltest.exe\",\"wmic.exe\",\"netstat.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where ParentProcessName endswith \"\\\\w3wp.exe\"\n| extend ChildBinary = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n| where ChildBinary in (_recon_or_shell)\n| extend CmdLower = tolower(CommandLine)\n| extend OilRigRecon = CmdLower has_any (\"whoami\",\"systeminfo\",\"net user\",\"net group\",\"net localgroup\",\"net view\",\"nltest /domain_trusts\",\"quser\",\"tasklist /v\",\"schtasks /create\",\"-encodedcommand\",\"-enc \",\"downloadstring\",\"iex\",\"invoke-expression\",\"invoke-webrequest\",\"certutil -urlcache\",\"bitsadmin /transfer\")\n| project TimeGenerated, Computer,\n          IISAppPoolUser = SubjectUserName,\n          ParentProcessName,\n          NewProcessName,\n          CommandLine,\n          OilRigRecon\n| order by TimeGenerated desc",
      "sigma_yaml": "title: APT34 (OilRig) IIS Web Shell \u2014 w3wp.exe Spawning Shell or Recon Binary\nid: 7c3a9e1f-4e6b-4a0e-9b6d-0f2a5d7e1c34\nstatus: experimental\ndescription: Detects Microsoft IIS worker process (w3wp.exe) spawning interactive shells or built-in discovery utilities \u2014 characteristic of APT34 / OilRig TwoFace, RGDoor and Karkoff web shells deployed against Internet-facing IIS / Exchange OWA / ECP.\nreferences:\n  - https://attack.mitre.org/groups/G0049/\n  - https://unit42.paloaltonetworks.com/unit42-twoface-webshell-persistent-access-point-lateral-movement/\n  - https://unit42.paloaltonetworks.com/unit42-rgdoor-isapi-backdoor-library/\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.persistence\n  - attack.t1505.003\n  - attack.execution\n  - attack.t1059.001\n  - attack.t1059.003\n  - attack.discovery\n  - attack.t1087\n  - attack.g0049\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent:\n    ParentImage|endswith: '\\w3wp.exe'\n  selection_child:\n    Image|endswith:\n      - '\\cmd.exe'\n      - '\\powershell.exe'\n      - '\\powershell_ise.exe'\n      - '\\pwsh.exe'\n      - '\\whoami.exe'\n      - '\\net.exe'\n      - '\\net1.exe'\n      - '\\systeminfo.exe'\n      - '\\quser.exe'\n      - '\\hostname.exe'\n      - '\\ipconfig.exe'\n      - '\\tasklist.exe'\n      - '\\schtasks.exe'\n      - '\\reg.exe'\n      - '\\bitsadmin.exe'\n      - '\\certutil.exe'\n      - '\\nltest.exe'\n      - '\\wmic.exe'\n      - '\\netstat.exe'\n  condition: selection_parent and selection_child\nfalsepositives:\n  - Legitimate IIS health-check / monitoring scripts that shell out from an app pool\n  - Web applications that intentionally run net.exe / whoami (rare, should be tightly scoped)\nlevel: high\n",
      "rationale": "Targeting Internet-facing IIS / Exchange via web shells is a defining APT34 tradecraft \u2014 Palo Alto Unit 42 documented TwoFace (2017) and RGDoor (2018) as OilRig-exclusive ASP.NET / IIS-module web shells, and Cisco Talos tied Karkoff (2019) to the same cluster. w3wp.exe spawning cmd/powershell/whoami is virtually never benign on production Exchange/IIS, and the recon-binary list mirrors the exact p"
    }
  ]
}