# intel/.llm_actor_uc_cache/58/58621875b14125e0f84d7437854b74257989357e.json
{
  "name": "Rhysida",
  "ucs": [
    {
      "title": "[LLM] Rhysida pre-ransomware NTDS.dit dump via ntdsutil IFM + AzCopy/MEGA staging",
      "description": "Hunts Rhysida's characteristic domain credential theft chain on a DC: ntdsutil 'ifm' snapshot of the Active Directory database followed within hours by exfiltration via AzCopy.exe or the MEGAsync/MEGAcmd client to actor-controlled cloud storage, immediately preceding ransomware detonation.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1003.003",
        "T1567.002",
        "T1059.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process) as parent_cmd from datamodel=Endpoint.Processes where (Processes.process_name=ntdsutil.exe Processes.process IN (\"*ac i ntds*\",\"*ifm*\",\"*create full*\")) OR (Processes.process_name IN (\"azcopy.exe\",\"MEGAsync.exe\",\"MEGAcmd*.exe\",\"megacmdshell.exe\",\"rclone.exe\") Processes.process IN (\"*copy*\",\"*sync*\",\"*--remote*\",\"*blob.core.windows.net*\",\"*mega.nz*\",\"*mega.co.nz*\")) by Processes.dest Processes.user Processes.process_name Processes.process Processes.parent_process_name _time | `drop_dm_object_name(Processes)` | eval bucket=floor(_time/14400) | stats values(process_name) as procs values(process) as cmds values(parent_process_name) as parents min(firstTime) as firstSeen max(lastTime) as lastSeen by dest user bucket | where mvcount(procs)>=2 AND match(mvjoin(procs,\",\"),\"ntdsutil\\.exe\") AND match(mvjoin(procs,\",\"),\"(azcopy|MEGA|rclone)\") | convert ctime(firstSeen) ctime(lastSeen) | sort - lastSeen",
      "kql": "let ntdsHosts = DeviceProcessEvents | where Timestamp > ago(7d) | where FileName =~ \"ntdsutil.exe\" or InitiatingProcessFileName =~ \"ntdsutil.exe\" | where ProcessCommandLine has_any (\"ifm\",\"create full\",\"ac i ntds\",\"activate instance ntds\") | project ntdsTime=Timestamp, DeviceId, DeviceName, ntdsCmd=ProcessCommandLine, ntdsUser=AccountName;\nlet exfilHosts = DeviceProcessEvents | where Timestamp > ago(7d) | where FileName in~ (\"azcopy.exe\",\"MEGAsync.exe\",\"MEGAcmd.exe\",\"MEGAcmdServer.exe\",\"megacmdshell.exe\",\"rclone.exe\") or ProcessCommandLine has_any (\"blob.core.windows.net\",\"mega.nz\",\"mega.co.nz\",\"--remote\",\"copy --recursive\") | project exfilTime=Timestamp, DeviceId, DeviceName, exfilProc=FileName, exfilCmd=ProcessCommandLine, exfilUser=AccountName;\nntdsHosts | join kind=inner exfilHosts on DeviceId | where exfilTime between (ntdsTime .. ntdsTime + 24h) | project ntdsTime, exfilTime, DeviceName, ntdsUser, exfilUser, ntdsCmd, exfilProc, exfilCmd | sort by exfilTime desc",
      "rationale": "CISA AA23-319A and Microsoft Threat Intelligence reporting on Rhysida (overlap with Vice Society) explicitly call out the ntdsutil 'IFM' snapshot tradecraft on domain controllers, followed by AzCopy or MEGA exfiltration before Rhysida.exe encryption. The ntdsutil+IFM keyword combo on a DC is an extremely high-signal anchor (legitimate admin use is rare and scheduled), and pairing it with AzCopy/ME"
    },
    {
      "title": "[LLM] Rhysida lateral deployment via PsExec + Defender tamper + scheduled-task ransomware payload",
      "description": "Detects Rhysida's signature ransomware deployment chain: an admin-context PsExec push from a beachhead server, immediately followed by PowerShell-based Microsoft Defender exclusion/disable commands and a scheduled task or remote service launching the ransomware binary (often masquerading as conhost/svchost) from C:\\Users\\Public, ProgramData or admin shares.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1021.002",
        "T1569.002",
        "T1562.001",
        "T1053.005",
        "T1486"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmds values(Processes.process_name) as procs values(Processes.parent_process_name) as parents from datamodel=Endpoint.Processes where Processes.process_name IN (\"PsExec.exe\",\"PsExec64.exe\",\"powershell.exe\",\"pwsh.exe\",\"schtasks.exe\",\"cmd.exe\") AND (Processes.process IN (\"*\\\\\\\\*-s -d -accepteula*\",\"*\\\\\\\\*cmd /c*\",\"*Set-MpPreference*DisableRealtimeMonitoring*\",\"*Add-MpPreference*ExclusionPath*\",\"*Set-MpPreference*-MAPSReporting 0*\",\"*schtasks*/create*/ru SYSTEM*\",\"*conhost.exe*-d*\",\"*ProgramData\\\\*.exe*\",\"*Users\\\\Public\\\\*.exe*\")) by Processes.dest Processes.user Processes.process_name Processes.process Processes.parent_process_name _time | `drop_dm_object_name(Processes)` | eval indicator=case(match(process,\"(?i)PsExec.{0,3}\\.exe\"),\"psexec\", match(process,\"(?i)DisableRealtimeMonitoring|ExclusionPath|MAPSReporting 0\"),\"defender_tamper\", match(process,\"(?i)schtasks.*/create.*/ru SYSTEM\"),\"sched_task\", match(process,\"(?i)(ProgramData|Users\\\\\\\\Public)\\\\\\\\.+\\\\.exe\"),\"suspicious_path\", true(),\"other\") | eval bucket=floor(_time/3600) | stats dc(indicator) as indicator_types values(indicator) as indicators values(process) as commandlines values(user) as users by dest bucket | where indicator_types>=3 AND mvfind(indicators,\"psexec\")>=0 | sort - bucket",
      "kql": "let window = 1h;\nlet psexec = DeviceProcessEvents | where Timestamp > ago(7d) | where FileName in~ (\"psexec.exe\",\"psexec64.exe\") or ProcessCommandLine has_any (\"-s -d -accepteula\",\"\\\\\\\\\") and InitiatingProcessFileName in~ (\"psexec.exe\",\"psexec64.exe\") | project t1=Timestamp, DeviceId, DeviceName, psexecCmd=ProcessCommandLine, psexecUser=AccountName;\nlet tamper = DeviceProcessEvents | where Timestamp > ago(7d) | where (FileName in~ (\"powershell.exe\",\"pwsh.exe\") and ProcessCommandLine has_any (\"DisableRealtimeMonitoring\",\"Add-MpPreference\",\"ExclusionPath\",\"ExclusionExtension\",\"MAPSReporting 0\",\"SubmitSamplesConsent 2\")) | project t2=Timestamp, DeviceId, tamperCmd=ProcessCommandLine;\nlet schedOrSvc = DeviceProcessEvents | where Timestamp > ago(7d) | where (FileName =~ \"schtasks.exe\" and ProcessCommandLine has \"/create\" and ProcessCommandLine has_any (\"/ru SYSTEM\",\"/ru \\\"NT AUTHORITY\\\\SYSTEM\\\"\")) or (FileName =~ \"sc.exe\" and ProcessCommandLine has_any (\"create\",\"start\")) | project t3=Timestamp, DeviceId, persistCmd=ProcessCommandLine;\nlet payloadDrop = DeviceFileEvents | where Timestamp > ago(7d) | where ActionType == \"FileCreated\" | where FolderPath has_any (\"\\\\ProgramData\\\\\",\"\\\\Users\\\\Public\\\\\",\"\\\\Windows\\\\Temp\\\\\",\"C$\\\\\") and FileName endswith \".exe\" | project t4=Timestamp, DeviceId, dropPath=FolderPath, dropFile=FileName, dropSha=SHA256;\npsexec\n| join kind=inner tamper on DeviceId | where t2 between (t1 .. t1+window)\n| join kind=inner schedOrSvc on DeviceId | where t3 between (t1 .. t1+window)\n| join kind=leftouter payloadDrop on DeviceId | where isnull(t4) or t4 between (t1 .. t1+window)\n| project t1, DeviceName, psexecUser, psexecCmd, tamperCmd, persistCmd, dropPath, dropFile, dropSha\n| sort by t1 desc",
      "rationale": "Rhysida (tracked alongside Vice Society lineage) is consistently documented by CISA, Sophos and SentinelOne as deploying its encryptor through PsExec from a staging host, disabling Defender via Set-MpPreference / Add-MpPreference exclusions, and launching the payload (often renamed conhost.exe or a 64-bit .exe in ProgramData/Public) via SCHTASKS /create /ru SYSTEM or sc.exe. Requiring all three be"
    }
  ]
}