# intel/.llm_actor_uc_cache/92/92bfec765dd9f4f726ca8c4c27750c40dd714367.json
{
  "name": "Embargo",
  "ucs": [
    {
      "title": "[LLM] Embargo / MDeployer safe-mode reboot to bypass EDR before encryption (BlackCat-lineage tradecraft)",
      "description": "Hunts Embargo's signature pre-encryption reboot-into-safe-mode chain (inherited from the BlackCat/ALPHV codebase) where the actor sets `safeboot` via bcdedit or registry, registers an auto-start RunOnce/SafeBoot entry pointing at MDeployer or the Rust encryptor, and forces a reboot \u2014 booting the host into Minimal/Network safe mode where Defender, MDE, and most EDR drivers do not load.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1562.009",
        "T1562.001",
        "T1490",
        "T1486",
        "T1543.003"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmd values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where (Processes.process_name=bcdedit.exe Processes.process IN (\"*safeboot*minimal*\",\"*safeboot*network*\")) OR (Processes.process_name=reg.exe Processes.process=\"*\\\\CurrentControlSet\\\\Control\\\\SafeBoot\\\\*\" Processes.process=\"*add*\") OR (Processes.process_name IN (powershell.exe,pwsh.exe) Processes.process IN (\"*Set-ItemProperty*SafeBoot*\",\"*New-ItemProperty*SafeBoot*\",\"*bcdedit*safeboot*\")) by host Processes.dest Processes.user Processes.process_name _time span=5m | `drop_dm_object_name(Processes)` | rename dest as host_dest | join type=inner host [ | tstats summariesonly=true count values(Processes.process) as shutdown_cmd from datamodel=Endpoint.Processes where Processes.process_name IN (shutdown.exe,cmd.exe,powershell.exe) Processes.process IN (\"*shutdown*/r*\",\"*shutdown*-r*\",\"*Restart-Computer*\") by host _time span=15m ] | where _time<=lastTime+900 AND _time>=firstTime-60 | eval signal=\"Embargo safe-mode reboot chain\" | table firstTime lastTime host user process_name cmd shutdown_cmd signal | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "// Embargo / MDeployer safe-mode reboot chain \u2014 bcdedit safeboot OR SafeBoot registry add, then forced reboot within 15 min on the same host\nlet Lookback = 7d;\nlet WindowMin = 15m;\nlet SafeBootStaging = DeviceProcessEvents\n    | where Timestamp > ago(Lookback)\n    | where AccountName !endswith \"$\"\n    | where (FileName =~ \"bcdedit.exe\"\n              and ProcessCommandLine has \"safeboot\"\n              and ProcessCommandLine has_any (\"minimal\", \"network\"))\n         or (FileName =~ \"reg.exe\"\n              and ProcessCommandLine has \"add\"\n              and ProcessCommandLine has @\"\\Control\\SafeBoot\\\")\n         or (FileName in~ (\"powershell.exe\", \"pwsh.exe\")\n              and ProcessCommandLine has_any (\"Set-ItemProperty\", \"New-ItemProperty\", \"bcdedit\")\n              and ProcessCommandLine has \"SafeBoot\")\n    | project SafeBootTime = Timestamp, DeviceId, DeviceName, AccountName,\n              StagingProcess = FileName,\n              StagingCmd = ProcessCommandLine,\n              StagingParent = InitiatingProcessFileName,\n              StagingParentCmd = InitiatingProcessCommandLine;\n// Reboot trigger on the SAME host shortly after\nDeviceProcessEvents\n| where Timestamp > ago(Lookback)\n| where AccountName !endswith \"$\"\n| where (FileName =~ \"shutdown.exe\" and ProcessCommandLine has_any (\"/r\", \"-r\"))\n     or (FileName in~ (\"powershell.exe\", \"pwsh.exe\") and ProcessCommandLine has \"Restart-Computer\")\n     or (FileName =~ \"cmd.exe\" and ProcessCommandLine has_all (\"shutdown\", \"/r\"))\n| join kind=inner SafeBootStaging on DeviceId\n| where Timestamp between (SafeBootTime .. SafeBootTime + WindowMin)\n| extend DelaySec = datetime_diff('second', Timestamp, SafeBootTime)\n| project SafeBootTime, RebootTime = Timestamp, DelaySec,\n          DeviceName, AccountName,\n          StagingProcess, StagingCmd, StagingParent, StagingParentCmd,\n          RebootProcess = FileName, RebootCmd = ProcessCommandLine\n| order by SafeBootTime desc",
      "sentinel_kql": "// Embargo / MDeployer safe-mode reboot chain \u2014 Sentinel SecurityEvent (4688) port\nlet Lookback = 7d;\nlet WindowMin = 15m;\nlet SafeBootStaging = SecurityEvent\n    | where TimeGenerated > ago(Lookback)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | where (NewProcessName endswith \"\\\\bcdedit.exe\"\n              and CommandLine has \"safeboot\"\n              and CommandLine has_any (\"minimal\", \"network\"))\n         or (NewProcessName endswith \"\\\\reg.exe\"\n              and CommandLine has \"add\"\n              and CommandLine has @\"\\Control\\SafeBoot\\\")\n         or (NewProcessName has_any (\"powershell.exe\", \"pwsh.exe\")\n              and CommandLine has_any (\"Set-ItemProperty\", \"New-ItemProperty\", \"bcdedit\")\n              and CommandLine has \"SafeBoot\")\n    | project SafeBootTime = TimeGenerated, Computer, Account,\n              StagingProcess = NewProcessName,\n              StagingCmd = CommandLine,\n              StagingParent = ParentProcessName;\nSecurityEvent\n| where TimeGenerated > ago(Lookback)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| where (NewProcessName endswith \"\\\\shutdown.exe\" and CommandLine has_any (\"/r\", \"-r\"))\n     or (NewProcessName has_any (\"powershell.exe\", \"pwsh.exe\") and CommandLine has \"Restart-Computer\")\n     or (NewProcessName endswith \"\\\\cmd.exe\" and CommandLine has_all (\"shutdown\", \"/r\"))\n| join kind=inner SafeBootStaging on Computer\n| where TimeGenerated between (SafeBootTime .. SafeBootTime + WindowMin)\n| extend DelaySec = datetime_diff('second', TimeGenerated, SafeBootTime)\n| project SafeBootTime, RebootTime = TimeGenerated, DelaySec,\n          Computer, Account,\n          StagingProcess, StagingCmd, StagingParent,\n          RebootProcess = NewProcessName, RebootCmd = CommandLine\n| order by SafeBootTime desc",
      "sigma_yaml": "title: Embargo Ransomware Safe-Mode Reboot Staging\nid: 8c5f7b2e-1d3a-4e2c-9b48-emba01safemode\nstatus: experimental\ndescription: Detects bcdedit / registry modification to enable Windows Safe Boot, a known precursor to Embargo (BlackCat-lineage) ransomware encryption to bypass EDR drivers.\nauthor: SOC Detection Eng\nreferences:\n  - https://www.welivesecurity.com/en/eset-research/embargo-ransomware-rock-n-rust/\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_bcdedit:\n    Image|endswith: '\\bcdedit.exe'\n    CommandLine|contains|all:\n      - 'safeboot'\n    CommandLine|contains:\n      - 'minimal'\n      - 'network'\n  selection_reg:\n    Image|endswith: '\\reg.exe'\n    CommandLine|contains|all:\n      - 'add'\n      - '\\Control\\SafeBoot\\'\n  selection_ps:\n    Image|endswith:\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n    CommandLine|contains: 'SafeBoot'\n  filter_user:\n    User|endswith: '$'\n  condition: (selection_bcdedit or selection_reg or selection_ps) and not filter_user\nfields:\n  - Image\n  - CommandLine\n  - ParentImage\n  - User\nfalsepositives:\n  - Legitimate troubleshooting by IT admins\n  - Windows Recovery / Imaging tools\nlevel: high\ntags:\n  - attack.defense_evasion\n  - attack.t1562.009\n  - attack.impact\n  - attack.t1486",
      "rationale": "Reboot-into-safe-mode is one of the few hard pivots that uniquely identifies Embargo and its BlackCat/ALPHV ancestor \u2014 REvil and Conti did not do it, and most modern crypters skip it. Embargo's MDeployer .NET loader specifically registers a SafeBoot RunOnce entry pointing at the encryptor, then issues `shutdown /r /t 0`. Pairing the two events on the same host within 15 minutes \u2014 both attributable",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] Embargo MS4Killer BYOVD \u2014 vulnerable kernel driver service install + EDR process termination",
      "description": "Hunts Embargo's MS4Killer (a.k.a. Probably-Killer / 'msdrv'-style) Bring-Your-Own-Vulnerable-Driver tooling: a kernel driver dropped to a non-system path (Temp / Public / ProgramData), registered as a service, then used (via IOCTL) to terminate Defender/EDR user-mode processes within minutes \u2014 the technique Embargo inherits from the BlackCat affiliate toolkit.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1562.001",
        "T1543.003",
        "T1068",
        "T1489",
        "T1027.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Filesystem.file_path) as file_path values(Filesystem.file_name) as file_name from datamodel=Endpoint.Filesystem where Filesystem.file_name=\"*.sys\" Filesystem.file_path IN (\"*\\\\Temp\\\\*\",\"*\\\\Users\\\\Public\\\\*\",\"*\\\\ProgramData\\\\*\",\"*\\\\AppData\\\\Local\\\\Temp\\\\*\") by host Filesystem.user _time span=10m | `drop_dm_object_name(Filesystem)` | join type=inner host [ | tstats summariesonly=true count values(Processes.process) as kill_cmds values(Processes.process_name) as killer_proc from datamodel=Endpoint.Processes where (Processes.process_name IN (taskkill.exe,sc.exe,net.exe,powershell.exe,pwsh.exe,wmic.exe) AND Processes.process IN (\"*MsMpEng*\",\"*MsSense*\",\"*SenseIR*\",\"*WinDefend*\",\"*Sense*\",\"*CSFalconService*\",\"*xagt*\",\"*SentinelAgent*\",\"*ekrn*\",\"*mfeesp*\",\"*MBAMService*\")) OR (Processes.parent_process_name=services.exe AND Processes.process IN (\"*\\\\Temp\\\\*.sys*\",\"*\\\\ProgramData\\\\*.sys*\",\"*\\\\Public\\\\*.sys*\")) by host _time span=30m ] | where lastTime>=firstTime AND (lastTime-firstTime)<1800 | eval signal=\"Embargo MS4Killer BYOVD candidate\" | table firstTime lastTime host user file_name file_path killer_proc kill_cmds signal | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "// Embargo MS4Killer BYOVD \u2014 non-system-path .sys dropped, kernel driver service installed, then EDR/AV processes killed on same host within 30 min\nlet Lookback = 7d;\nlet WindowMin = 30m;\nlet EdrTargets = dynamic([\n    \"MsMpEng.exe\", \"MsSense.exe\", \"SenseIR.exe\", \"SenseCncProxy.exe\", \"WinDefend\",\n    \"CSFalconService.exe\", \"CSFalconContainer.exe\", \"xagt.exe\", \"SentinelAgent.exe\",\n    \"SentinelServiceHost.exe\", \"ekrn.exe\", \"egui.exe\", \"mfeesp.exe\", \"masvc.exe\",\n    \"MBAMService.exe\", \"AvastSvc.exe\", \"avp.exe\", \"sophosfs.exe\", \"SAVService.exe\",\n    \"TmCCSF.exe\", \"PccNTMon.exe\"\n  ]);\nlet KnownVulnerableDrivers = dynamic([\n    \"truesight.sys\", \"gmer.sys\", \"viragt64.sys\", \"procexp152.sys\", \"procexp.sys\",\n    \"rentdrv2.sys\", \"dbutil_2_3.sys\", \"dbutildrv2.sys\", \"ene.sys\", \"asrdrv101.sys\",\n    \"asrdrv102.sys\", \"asrdrv103.sys\", \"WnBios.sys\", \"WnBios64.sys\", \"mhyprot2.sys\",\n    \"daxin.sys\", \"kavg.sys\", \"ms4killer.sys\"\n  ]);\nlet DriverDrop = DeviceFileEvents\n    | where Timestamp > ago(Lookback)\n    | where ActionType in (\"FileCreated\", \"FileRenamed\")\n    | where FileName endswith \".sys\"\n    | where FolderPath has_any (@\"\\Temp\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Users\\Public\\\",\n                                 @\"\\ProgramData\\\", @\"\\Downloads\\\", @\"\\PerfLogs\\\")\n        or FileName in~ (KnownVulnerableDrivers)\n    | where InitiatingProcessAccountName !endswith \"$\"\n    | project DropTime = Timestamp, DeviceId, DeviceName,\n              DroppedDriver = FileName, DropPath = FolderPath, DropSHA256 = SHA256,\n              DropperProcess = InitiatingProcessFileName,\n              DropperCmd = InitiatingProcessCommandLine,\n              DropperUser = InitiatingProcessAccountName;\n// Service registration for that driver (HKLM\\SYSTEM\\CurrentControlSet\\Services\\<name>\\ImagePath)\nlet ServiceReg = DeviceRegistryEvents\n    | where Timestamp > ago(Lookback)\n    | where ActionType in (\"RegistryValueSet\", \"RegistryKeyCreated\")\n    | where RegistryKey has @\"\\CurrentControlSet\\Services\\\"\n    | where RegistryValueName =~ \"ImagePath\"\n    | where RegistryValueData has_any (@\"\\Temp\\\", @\"\\Public\\\", @\"\\ProgramData\\\",\n                                        @\"\\AppData\\Local\\Temp\\\", @\"\\Downloads\\\")\n    | where RegistryValueData has \".sys\"\n    | project SvcRegTime = Timestamp, DeviceId, DeviceName,\n              SvcKey = RegistryKey, SvcImagePath = tostring(RegistryValueData),\n              SvcRegisterer = InitiatingProcessFileName,\n              SvcRegistererCmd = InitiatingProcessCommandLine;\n// EDR/AV process termination on same host\nlet EdrKill = DeviceProcessEvents\n    | where Timestamp > ago(Lookback)\n    | where (FileName in~ (\"taskkill.exe\", \"sc.exe\", \"net.exe\", \"wmic.exe\",\n                            \"powershell.exe\", \"pwsh.exe\")\n              and (ProcessCommandLine has_any (EdrTargets)\n                   or ProcessCommandLine has_any (\"WinDefend\", \"Sense\", \"MDE\", \"Defender\")))\n    | project KillTime = Timestamp, DeviceId, DeviceName, AccountName,\n              KillProcess = FileName, KillCmd = ProcessCommandLine,\n              KillParent = InitiatingProcessFileName;\nDriverDrop\n| join kind=inner ServiceReg on DeviceId\n| where SvcRegTime between (DropTime .. DropTime + WindowMin)\n| join kind=inner EdrKill on DeviceId\n| where KillTime between (DropTime .. DropTime + WindowMin)\n| project DropTime, SvcRegTime, KillTime,\n          DeviceName, DropperUser,\n          DroppedDriver, DropPath, DropSHA256,\n          DropperProcess, DropperCmd,\n          SvcImagePath, SvcRegisterer,\n          KillProcess, KillCmd\n| order by DropTime desc",
      "sentinel_kql": "// Embargo MS4Killer BYOVD \u2014 Sentinel ASIM port (ImFileEvent + ImRegistryEvent + ImProcessCreate)\nlet Lookback = 7d;\nlet WindowMin = 30m;\nlet EdrTargets = dynamic([\n    \"MsMpEng\", \"MsSense\", \"SenseIR\", \"SenseCncProxy\", \"WinDefend\",\n    \"CSFalconService\", \"xagt\", \"SentinelAgent\", \"ekrn\", \"mfeesp\",\n    \"MBAMService\", \"AvastSvc\", \"avp\", \"sophosfs\", \"TmCCSF\"\n  ]);\nlet KnownVulnerableDrivers = dynamic([\n    \"truesight.sys\", \"gmer.sys\", \"viragt64.sys\", \"procexp152.sys\",\n    \"rentdrv2.sys\", \"dbutil_2_3.sys\", \"dbutildrv2.sys\", \"ene.sys\",\n    \"asrdrv101.sys\", \"WnBios.sys\", \"mhyprot2.sys\", \"ms4killer.sys\"\n  ]);\nlet DriverDrop = ImFileEvent\n    | where TimeGenerated > ago(Lookback)\n    | where EventType in (\"FileCreated\", \"FileRenamed\")\n    | where TargetFileName endswith \".sys\"\n    | where TargetFilePath has_any (@\"\\Temp\\\", @\"\\AppData\\Local\\Temp\\\",\n                                     @\"\\Users\\Public\\\", @\"\\ProgramData\\\", @\"\\Downloads\\\")\n        or TargetFileName in~ (KnownVulnerableDrivers)\n    | where ActorUsername !endswith \"$\"\n    | project DropTime = TimeGenerated, DvcHostname,\n              DroppedDriver = TargetFileName, DropPath = TargetFilePath,\n              DropSHA256 = TargetFileSHA256,\n              DropperProcess = ActingProcessName,\n              DropperCmd = ActingProcessCommandLine,\n              DropperUser = ActorUsername;\nlet ServiceReg = ImRegistryEvent\n    | where TimeGenerated > ago(Lookback)\n    | where RegistryKey has @\"\\CurrentControlSet\\Services\\\"\n    | where RegistryValue =~ \"ImagePath\"\n    | where tostring(RegistryValueData) has_any (@\"\\Temp\\\", @\"\\Public\\\",\n                                                  @\"\\ProgramData\\\", @\"\\Downloads\\\")\n    | where tostring(RegistryValueData) has \".sys\"\n    | project SvcRegTime = TimeGenerated, DvcHostname,\n              SvcKey = RegistryKey,\n              SvcImagePath = tostring(RegistryValueData),\n              SvcRegisterer = ActingProcessName;\nlet EdrKill = ImProcessCreate\n    | where TimeGenerated > ago(Lookback)\n    | where TargetProcessName has_any (\"taskkill.exe\", \"sc.exe\", \"net.exe\", \"wmic.exe\",\n                                         \"powershell.exe\", \"pwsh.exe\")\n    | where TargetProcessCommandLine has_any (EdrTargets)\n         or TargetProcessCommandLine has_any (\"WinDefend\", \"Sense\", \"Defender\")\n    | project KillTime = TimeGenerated, DvcHostname, ActorUsername,\n              KillProcess = TargetProcessName,\n              KillCmd = TargetProcessCommandLine;\nDriverDrop\n| join kind=inner ServiceReg on DvcHostname\n| where SvcRegTime between (DropTime .. DropTime + WindowMin)\n| join kind=inner EdrKill on DvcHostname\n| where KillTime between (DropTime .. DropTime + WindowMin)\n| project DropTime, SvcRegTime, KillTime,\n          DvcHostname, DropperUser,\n          DroppedDriver, DropPath, DropSHA256,\n          DropperProcess, DropperCmd,\n          SvcImagePath, SvcRegisterer,\n          KillProcess, KillCmd\n| order by DropTime desc",
      "sigma_yaml": "",
      "rationale": "Embargo's distinguishing tooling is the Rust-based MS4Killer (a.k.a. 'Probably-Killer' / 'MS4'), which ESET and TrendMicro have publicly attributed to the group as inherited from the BlackCat/ALPHV affiliate kit. Its operational signature \u2014 drop a vulnerable signed driver (often impersonating Defender by name 'MS4Killer.sys'), register it via Service Control Manager from a non-system path, then ca"
    }
  ]
}