# intel/.llm_actor_uc_cache/ca/caec9fb5e1499fb42747e15f47e0eeb55b01696f.json
{
  "name": "Windshift",
  "ucs": [
    {
      "title": "[LLM] Windshift / Bahamut VBScript spearphishing stager \u2192 HKCU Run persistence chain",
      "description": "Hunts the canonical Windshift (a.k.a. Bahamut, MITRE G0112) intrusion chain against Middle-East targets: a spearphishing attachment/link executed from Outlook/archive utilities spawns wscript.exe or cscript.exe (their preferred VBScript stager \u2014 T1059.005), which then writes a HKCU\\...\\Run value pointing into %APPDATA%, %ProgramData% or %PUBLIC% within minutes. The wscript\u2192Run-key correlation is what differentiates Windshift from generic VBS noise.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1566.002",
        "T1204.002",
        "T1059.005",
        "T1547.001",
        "T1027"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as stager_time from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"outlook.exe\",\"winrar.exe\",\"7zfm.exe\",\"explorer.exe\",\"chrome.exe\",\"msedge.exe\") Processes.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\") NOT Processes.user IN (\"*$\",\"SYSTEM\",\"LOCAL SERVICE\",\"NETWORK SERVICE\") by Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process \n| `drop_dm_object_name(\"Processes\")` \n| join type=inner dest [\n    | tstats summariesonly=t count min(_time) as persistence_time values(Registry.registry_value_data) as reg_data values(Registry.registry_path) as reg_path from datamodel=Endpoint.Registry where Registry.registry_path IN (\"*\\\\CurrentVersion\\\\Run\\\\*\",\"*\\\\CurrentVersion\\\\RunOnce\\\\*\") (Registry.registry_value_data=\"*\\\\AppData\\\\*\" OR Registry.registry_value_data=\"*\\\\ProgramData\\\\*\" OR Registry.registry_value_data=\"*\\\\Public\\\\*\" OR Registry.registry_value_data=\"*.vbs*\" OR Registry.registry_value_data=\"*.js*\" OR Registry.registry_value_data=\"*wscript*\" OR Registry.registry_value_data=\"*cscript*\") by Registry.dest \n    | `drop_dm_object_name(\"Registry\")`\n  ] \n| eval delay_minutes = round((persistence_time - stager_time)/60,1) \n| where delay_minutes >= 0 AND delay_minutes <= 10 \n| table stager_time persistence_time delay_minutes dest user parent_process_name process_name process reg_path reg_data \n| sort - stager_time",
      "kql": "// Windshift / Bahamut (G0112) \u2014 spearphish-spawned VBS stager + HKCU Run persistence within 10 min\nlet LookbackDays = 14d;\nlet WindowMinutes = 10;\nlet _mail_or_archive_parents = dynamic([\"outlook.exe\",\"winrar.exe\",\"7zfm.exe\",\"winzip32.exe\",\"explorer.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\"]);\nlet _vb_stagers = dynamic([\"wscript.exe\",\"cscript.exe\",\"mshta.exe\"]);\nlet StagerExec =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where InitiatingProcessFileName in~ (_mail_or_archive_parents)\n    | where FileName in~ (_vb_stagers)\n    | where AccountName !endswith \"$\"                       // skip machine accts\n    | where AccountName !in~ (\"system\",\"local service\",\"network service\")\n    | project StagerTime = Timestamp, DeviceId, DeviceName, AccountName,\n              ParentName = InitiatingProcessFileName,\n              StagerName = FileName,\n              StagerCmd  = ProcessCommandLine,\n              StagerSHA256 = SHA256;\nDeviceRegistryEvents\n| where Timestamp > ago(LookbackDays)\n| where ActionType == \"RegistryValueSet\"\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\")\n| where RegistryValueData has_any (@\"\\AppData\\\", @\"\\ProgramData\\\", @\"\\Public\\\",\n                                   \".vbs\",\".vbe\",\".js\",\".jse\",\".lnk\",\n                                   \"wscript\",\"cscript\",\"mshta\")\n| where InitiatingProcessFileName !in~ (\"msiexec.exe\",\"explorer.exe\",\"setup.exe\")  // legit installers\n| join kind=inner StagerExec on DeviceId\n| where Timestamp between (StagerTime .. StagerTime + WindowMinutes * 1m)\n| project StagerTime, PersistenceTime = Timestamp,\n          DelayMin = datetime_diff('minute', Timestamp, StagerTime),\n          DeviceName, AccountName, ParentName, StagerName, StagerCmd, StagerSHA256,\n          RegistryKey, RegistryValueName, RegistryValueData,\n          PersistenceWriter = InitiatingProcessFileName\n| order by StagerTime desc",
      "sentinel_kql": "// Windshift / Bahamut (G0112) \u2014 spearphish-spawned VBS stager + HKCU Run persistence within 10 min\nlet LookbackDays = 14d;\nlet WindowMinutes = 10;\nlet _vb_stagers = dynamic([\"wscript.exe\",\"cscript.exe\",\"mshta.exe\"]);\nlet _mail_or_archive_parents = dynamic([\"outlook.exe\",\"winrar.exe\",\"7zfm.exe\",\"winzip32.exe\",\"explorer.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\"]);\nlet StagerExec =\n    SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ParentLeaf = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1])),\n             ChildLeaf  = tolower(tostring(split(NewProcessName,    \"\\\\\")[-1]))\n    | where ParentLeaf in (_mail_or_archive_parents)\n    | where ChildLeaf  in (_vb_stagers)\n    | project StagerTime = TimeGenerated, Computer, Account,\n              ParentName = ParentLeaf, StagerName = ChildLeaf,\n              StagerCmd  = CommandLine;\nImRegistryEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventType == \"RegistryValueSet\"\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\")\n| where RegistryValueData has_any (@\"\\AppData\\\", @\"\\ProgramData\\\", @\"\\Public\\\",\n                                   \".vbs\",\".vbe\",\".js\",\".jse\",\".lnk\",\n                                   \"wscript\",\"cscript\",\"mshta\")\n| where ActingProcessName !in~ (\"msiexec.exe\",\"explorer.exe\",\"setup.exe\")\n| join kind=inner StagerExec on $left.DvcHostname == $right.Computer\n| where TimeGenerated between (StagerTime .. StagerTime + WindowMinutes * 1m)\n| project StagerTime, PersistenceTime = TimeGenerated,\n          DelayMin = datetime_diff('minute', TimeGenerated, StagerTime),\n          Computer, Account, ParentName, StagerName, StagerCmd,\n          RegistryKey, RegistryValue, RegistryValueData,\n          PersistenceWriter = ActingProcessName\n| order by StagerTime desc",
      "sigma_yaml": "",
      "rationale": "Windshift/Bahamut's published intrusion chains (SANS 2018, Objective-See WindTail 2018-2019, Talos Bahamut 2022) consistently show: (1) Middle-East-targeted spearphish (T1566.001/.002) using VBS-bearing attachments and .lnk wrappers; (2) wscript/cscript spawned from Outlook or archive tools (T1059.005, T1204.002); (3) immediate persistence via HKCU\\...\\Run pointing at user-writable directories (T1"
    },
    {
      "title": "[LLM] Windshift post-exploit reconnaissance \u2014 VBS-spawned WMI/native security-product enumeration",
      "description": "Hunts Windshift/Bahamut's distinctive on-target triage step: after the VBS stager lands, the actor's implant uses wmic.exe / tasklist / native binaries to enumerate AV products, running processes and OS info (T1518.001, T1057, T1082, T1047) before deciding to drop the second-stage. Catches the recon burst spawned from a script interpreter that itself was launched from a mail/archive parent \u2014 a chain unique to user-execution-driven actors like G0112.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "medium",
      "techs": [
        "T1047",
        "T1057",
        "T1082",
        "T1518",
        "T1518.001",
        "T1059.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count values(Processes.process) as recon_cmds dc(Processes.process_name) as distinct_recon_bins from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\") (Processes.process_name IN (\"wmic.exe\",\"tasklist.exe\",\"systeminfo.exe\",\"net.exe\",\"netstat.exe\",\"reg.exe\",\"whoami.exe\") OR (Processes.process_name=\"powershell.exe\" AND (Processes.process=\"*Win32_Product*\" OR Processes.process=\"*AntiVirusProduct*\" OR Processes.process=\"*Get-Process*\"))) NOT Processes.user IN (\"*$\",\"SYSTEM\") by Processes.dest Processes.user Processes.parent_process_name _time span=10m \n| `drop_dm_object_name(\"Processes\")` \n| where distinct_recon_bins >= 3 \n| sort - _time",
      "kql": "// Windshift / Bahamut (G0112) \u2014 VBS-stager-driven host triage burst (>=3 distinct recon bins in 10 min)\nlet Lookback = 14d;\nlet RunBin = 10m;\nlet _recon_bins = dynamic([\"wmic.exe\",\"tasklist.exe\",\"systeminfo.exe\",\"net.exe\",\"net1.exe\",\"netstat.exe\",\"reg.exe\",\"whoami.exe\",\"hostname.exe\",\"ipconfig.exe\",\"nltest.exe\"]);\nlet _av_keywords = dynamic([\"AntiVirusProduct\",\"AntiSpywareProduct\",\"FirewallProduct\",\"securitycenter2\",\"Win32_Product\",\"Get-MpComputerStatus\",\"Get-MpPreference\"]);\nDeviceProcessEvents\n| where Timestamp > ago(Lookback)\n| where AccountName !endswith \"$\"\n| where InitiatingProcessFileName in~ (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\")\n| where FileName in~ (_recon_bins)\n   or (FileName in~ (\"powershell.exe\",\"pwsh.exe\")\n       and ProcessCommandLine has_any (_av_keywords))\n   or (FileName =~ \"wmic.exe\"\n       and ProcessCommandLine has_any (\"antivirus\",\"antispyware\",\"firewallproduct\",\n                                       \"process list\",\"os get\",\"computersystem\",\n                                       \"product get\"))\n| summarize ReconCount       = count(),\n            DistinctReconBins = dcount(FileName),\n            ReconBins         = make_set(FileName, 25),\n            SampleCmds        = make_set(ProcessCommandLine, 25),\n            FirstSeen         = min(Timestamp),\n            LastSeen          = max(Timestamp)\n            by DeviceName, AccountName, InitiatingProcessFileName,\n               InitiatingProcessCommandLine, bin(Timestamp, RunBin)\n| where DistinctReconBins >= 3                               // burst behaviour, not single tool\n| order by FirstSeen desc",
      "sentinel_kql": "// Windshift / Bahamut (G0112) \u2014 VBS-stager-driven host triage burst\nlet Lookback = 14d;\nlet RunBin = 10m;\nlet _recon_bins = dynamic([\"wmic.exe\",\"tasklist.exe\",\"systeminfo.exe\",\"net.exe\",\"net1.exe\",\"netstat.exe\",\"reg.exe\",\"whoami.exe\",\"hostname.exe\",\"ipconfig.exe\",\"nltest.exe\"]);\nlet _av_keywords = dynamic([\"AntiVirusProduct\",\"AntiSpywareProduct\",\"FirewallProduct\",\"securitycenter2\",\"Win32_Product\",\"Get-MpComputerStatus\",\"Get-MpPreference\"]);\nSecurityEvent\n| where TimeGenerated > ago(Lookback)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| extend ParentLeaf = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1])),\n         ChildLeaf  = tolower(tostring(split(NewProcessName,    \"\\\\\")[-1]))\n| where ParentLeaf in~ (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\")\n| where ChildLeaf in~ (_recon_bins)\n   or (ChildLeaf in~ (\"powershell.exe\",\"pwsh.exe\")\n       and CommandLine has_any (_av_keywords))\n   or (ChildLeaf =~ \"wmic.exe\"\n       and CommandLine has_any (\"antivirus\",\"antispyware\",\"firewallproduct\",\n                                \"process list\",\"os get\",\"computersystem\",\n                                \"product get\"))\n| summarize ReconCount       = count(),\n            DistinctReconBins = dcount(ChildLeaf),\n            ReconBins         = make_set(ChildLeaf, 25),\n            SampleCmds        = make_set(CommandLine, 25),\n            FirstSeen         = min(TimeGenerated),\n            LastSeen          = max(TimeGenerated)\n            by Computer, Account, ParentLeaf, bin(TimeGenerated, RunBin)\n| where DistinctReconBins >= 3\n| order by FirstSeen desc",
      "sigma_yaml": "",
      "rationale": "Windshift implants (WindShift Win backdoor, Bahamut Windows variants documented by Cisco Talos and SANS) profile the host before second-stage delivery \u2014 calling out specifically to WMI for AV/firewall product enumeration (T1518.001) and process listings (T1057). Generic recon binaries on their own are noise, but a *burst* of \u22653 distinct recon binaries inside 10 minutes spawned directly by wscript/"
    }
  ]
}