# intel/.llm_actor_uc_cache/ca/caec9fb5e1499fb42747e15f47e0eeb55b01696f.json
{
  "name": "Windshift",
  "ucs": [
    {
      "title": "[LLM] Windshift / Bahamut (G0112) WindTail macOS dropper: quarantined .app spawned from Archive Utility with outbound HTTP beacon",
      "description": "Hunts the WindTail (a.k.a. WINDTAIL/Bahamut) macOS backdoor delivery chain unique to Windshift in the Middle East: a phishing link drops a ZIP whose embedded .app is auto-decompressed by Archive Utility and launched via a registered custom URL scheme, then beacons over HTTP. This combination of macOS targeting + Archive-Utility-spawned ad-hoc-signed bundle is a Windshift-specific signature, distinguishing it from generic Mac adware.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.002",
        "T1204.002",
        "T1036.001",
        "T1071.001",
        "T1547.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.process_path) as proc_path from datamodel=Endpoint.Processes where Processes.os=\"macOS\" AND Processes.parent_process_name IN (\"Archive Utility\",\"BOMArchiveHelper\",\"Safari\",\"Mail\",\"Finder\",\"Google Chrome\") AND (Processes.process_path=\"*/Users/*/Downloads/*.app/Contents/MacOS/*\" OR Processes.process_path=\"*/private/tmp/*.app/Contents/MacOS/*\" OR Processes.process_path=\"*/Volumes/*.app/Contents/MacOS/*\") AND (Processes.signature=\"\" OR Processes.signature=\"adhoc\" OR Processes.signature_verified=false) by host Processes.user Processes.process_name Processes.parent_process_name Processes.dest | `drop_dm_object_name(Processes)` | join type=inner host [| tstats summariesonly=true values(All_Traffic.dest) as c2 values(All_Traffic.url) as urls from datamodel=Network_Traffic where All_Traffic.app=\"curl\" OR All_Traffic.app=\"python\" OR (All_Traffic.dest_port IN (80,443) AND All_Traffic.bytes_out<2048 AND All_Traffic.url=\"*.php*\") by host | rename All_Traffic.* as *] | where lastTime-firstTime < 600 | table firstTime host user process_name parent_process_name proc_path c2 urls",
      "kql": "let suspectApps = DeviceProcessEvents | where Timestamp > ago(14d) | where DeviceType =~ \"Mac\" or InitiatingProcessVersionInfoCompanyName has \"Apple\" | where FolderPath matches regex @\"\\.app/Contents/MacOS/\" and (FolderPath startswith \"/Users/\" or FolderPath startswith \"/private/tmp/\" or FolderPath startswith \"/Volumes/\") | where InitiatingProcessFileName in~ (\"Archive Utility\",\"BOMArchiveHelper\",\"Safari\",\"Mail\",\"Finder\",\"Google Chrome\",\"firefox\",\"Microsoft Edge\") | where isempty(ProcessVersionInfoCompanyName) or ProcessVersionInfoCompanyName !has \"Apple\" | project Timestamp, DeviceId, DeviceName, AccountName, FileName, FolderPath, ProcessCommandLine, InitiatingProcessFileName, InitiatingProcessId; let beacons = DeviceNetworkEvents | where Timestamp > ago(14d) | where RemoteUrl has_any (\".php\",\"/checkupdate\",\"/get/\",\"/string/\",\"/usr/\") or (RemotePort in (80,443) and InitiatingProcessFileName in~ (\"curl\",\"python\",\"python3\")) | project NetTime=Timestamp, DeviceId, RemoteUrl, RemoteIP, InitiatingProcessId; suspectApps | join kind=inner beacons on DeviceId, $left.InitiatingProcessId == $right.InitiatingProcessId | where datetime_diff('minute', NetTime, Timestamp) between (0 .. 30) | project Timestamp, DeviceName, AccountName, FolderPath, ProcessCommandLine, RemoteUrl, RemoteIP",
      "rationale": "Windshift is one of the very few APTs that has consistently weaponized native macOS backdoors (WindTail/WindTape) against Middle Eastern targets. Objective-See's analyses (Wardle, Dec 2018 / Jan 2019) document the exact chain modeled here: ZIP delivered via spearphishing link, decompressed by Archive Utility, executed via custom URL scheme handler, ad-hoc signed bundle stored under ~/Downloads, th"
    },
    {
      "title": "[LLM] Windshift / Bahamut spearphish-to-persistence: VBScript launched by Office/browser writes Run-key auto-start",
      "description": "Detects the Windows side of the Windshift kill chain documented by SANS (Aug 2018): spearphishing email or link drops a VBScript/JScript stager that is executed via wscript.exe from a user-writable temp path, then registers a Run / RunOnce key pointing back at the same script for surveillance persistence. The temporal correlation of Office/browser \u2192 wscript from Temp \u2192 Run-key write within minutes is the actor-specific signature.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1566.002",
        "T1204.002",
        "T1059.005",
        "T1547.001",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true min(_time) as procTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent from datamodel=Endpoint.Processes where Processes.process_name IN (\"wscript.exe\",\"cscript.exe\") AND Processes.parent_process_name IN (\"outlook.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"explorer.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"iexplore.exe\") AND (Processes.process=\"*.vbs*\" OR Processes.process=\"*.vbe*\" OR Processes.process=\"*.js*\" OR Processes.process=\"*.jse*\") AND (Processes.process=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Processes.process=\"*\\\\AppData\\\\Roaming\\\\*\" OR Processes.process=\"*\\\\Downloads\\\\*\" OR Processes.process=\"*\\\\ProgramData\\\\*\" OR Processes.process=\"*\\\\Users\\\\Public\\\\*\") by host Processes.user Processes.process_id | `drop_dm_object_name(Processes)` | join type=inner host [| tstats summariesonly=true min(_time) as regTime values(Registry.registry_value_data) as run_value values(Registry.registry_key_name) as run_key from datamodel=Endpoint.Registry where Registry.action=\"modified\" AND Registry.registry_path=\"*\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run*\" AND (Registry.registry_value_data=\"*.vbs*\" OR Registry.registry_value_data=\"*.vbe*\" OR Registry.registry_value_data=\"*wscript*\" OR Registry.registry_value_data=\"*cscript*\") by host | rename Registry.* as *] | eval delta=abs(regTime-procTime) | where delta < 1800 | table procTime regTime delta host user parent cmdline run_key run_value",
      "kql": "let stager = DeviceProcessEvents | where Timestamp > ago(7d) | where FileName in~ (\"wscript.exe\",\"cscript.exe\") | where InitiatingProcessFileName in~ (\"outlook.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"explorer.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"iexplore.exe\") | where ProcessCommandLine has_any (\".vbs\",\".vbe\",\".js\",\".jse\") | where ProcessCommandLine matches regex @\"(?i)\\\\(AppData\\\\Local\\\\Temp|AppData\\\\Roaming|Downloads|ProgramData|Users\\\\Public)\\\\\" | project ProcTime=Timestamp, DeviceId, DeviceName, AccountName, ProcessCommandLine, InitiatingProcessFileName, ProcessId; let persist = DeviceRegistryEvents | where Timestamp > ago(7d) | where ActionType == \"RegistryValueSet\" | where RegistryKey has_any (@\"\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\", @\"\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\") | where RegistryValueData has_any (\".vbs\",\".vbe\",\"wscript\",\"cscript\") | project RegTime=Timestamp, DeviceId, RegistryKey, RegistryValueName, RegistryValueData; stager | join kind=inner persist on DeviceId | where datetime_diff('minute', RegTime, ProcTime) between (-30 .. 30) | join kind=leftouter (EmailEvents | where Timestamp > ago(14d) | where DeliveryAction == \"Delivered\" | project EmailTime=Timestamp, RecipientEmailAddress, SenderFromAddress, Subject, NetworkMessageId) on $left.AccountName == $right.RecipientEmailAddress | project ProcTime, DeviceName, AccountName, InitiatingProcessFileName, ProcessCommandLine, RegistryValueData, SenderFromAddress, Subject",
      "rationale": "Windshift's Windows surveillance chain documented by SANS ICS (Karim, Aug 2018) and follow-on Bahamut reporting hinges on simple but distinctive components: VBScript/JScript stagers under user-writable directories, executed via the standard wscript host, with persistence via classic Run/RunOnce keys \u2014 not scheduled tasks, services, or WMI subscriptions favored by louder crews. Correlating an Offic"
    }
  ]
}