# intel/.llm_actor_uc_cache/a6/a623b6e8bce5bbfd8ad62200ffe87cf54110c81b.json
{
  "name": "Gorgon Group",
  "ucs": [
    {
      "title": "[LLM] Gorgon Group: Office macro \u2192 hidden-window scripting host \u2192 process hollowing into RegAsm/MSBuild (.NET LOLBin)",
      "description": "Hunts the canonical Gorgon Group (Subaat) infection chain reported by Unit 42 in 2018 and observed continuously since: a weaponized Office document drops a VBScript/PowerShell stub which downloads a commodity RAT (NJRat, RevengeRAT, LokiBot, Remcos, AZORult) and process-hollows it into a signed .NET utility \u2014 most often RegAsm.exe, RegSvcs.exe, MSBuild.exe, InstallUtil.exe or aspnet_compiler.exe. Catches the chain by joining Stage 1 (Office \u2192 hidden script host with download/IEX/base64 markers) to Stage 2 (script host spawning a .NET LOLBin, the hollow target) within a 5-minute window on the s",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.001",
        "T1059.005",
        "T1564.003",
        "T1140",
        "T1105",
        "T1055.012"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true min(_time) as stage1_time count from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"eqnedt32.exe\") AND Processes.process_name IN (\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\") AND (Processes.process=\"*hidden*\" OR Processes.process=\"*-w 1*\" OR Processes.process=\"*-windowstyle*\" OR Processes.process=\"*FromBase64String*\" OR Processes.process=\"*DownloadString*\" OR Processes.process=\"*DownloadFile*\" OR Processes.process=\"*Invoke-Expression*\" OR Processes.process=\"*IEX*\" OR Processes.process=\"*-EncodedCommand*\" OR Processes.process=\"*pastebin*\" OR Processes.process=\"*bit.ly*\") by Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process | `drop_dm_object_name(Processes)` | rename parent_process_name as office_parent, process_name as script_host, process as script_cmd | join type=inner dest [ | tstats summariesonly=true min(_time) as hollow_time from datamodel=Endpoint.Processes where Processes.process_name IN (\"regasm.exe\",\"regsvcs.exe\",\"msbuild.exe\",\"installutil.exe\",\"aspnet_compiler.exe\",\"caspol.exe\") AND Processes.parent_process_name IN (\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\") by Processes.dest Processes.process_name Processes.process Processes.parent_process_name | `drop_dm_object_name(Processes)` | rename process_name as hollow_target, process as hollow_cmd, parent_process_name as hollow_parent ] | where hollow_time >= stage1_time AND (hollow_time - stage1_time) <= 300 | table dest user stage1_time hollow_time office_parent script_host script_cmd hollow_target hollow_parent hollow_cmd",
      "kql": "let LookbackDays = 7d;\nlet WindowMin = 5m;\nlet OfficeApps     = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"eqnedt32.exe\"]);\nlet ScriptHosts    = dynamic([\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\"]);\nlet HollowTargets  = dynamic([\"regasm.exe\",\"regsvcs.exe\",\"msbuild.exe\",\"installutil.exe\",\"aspnet_compiler.exe\",\"caspol.exe\"]);\nlet Stage1 = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where InitiatingProcessFileName in~ (OfficeApps)\n    | where FileName in~ (ScriptHosts)\n    | where ProcessCommandLine has_any (\"hidden\",\"-w 1\",\"-windowstyle\",\"FromBase64String\",\"DownloadString\",\"DownloadFile\",\"Invoke-Expression\",\" iex \",\"-EncodedCommand\",\"pastebin\",\"bit.ly\",\"j.mp\")\n    | project Stage1Time = Timestamp, DeviceId, DeviceName, AccountName,\n              OfficeParent = InitiatingProcessFileName,\n              ScriptHost   = FileName,\n              ScriptCmd    = ProcessCommandLine,\n              ScriptPid    = ProcessId;\nDeviceProcessEvents\n| where Timestamp > ago(LookbackDays)\n| where FileName in~ (HollowTargets)\n| where InitiatingProcessFileName in~ (ScriptHosts)\n| join kind=inner Stage1 on DeviceId\n| where Timestamp between (Stage1Time .. Stage1Time + WindowMin)\n| project DeviceName, AccountName, Stage1Time, HollowTime = Timestamp,\n          DelaySec    = datetime_diff('second', Timestamp, Stage1Time),\n          OfficeParent, ScriptHost, ScriptCmd,\n          HollowParent = InitiatingProcessFileName,\n          HollowTarget = FileName,\n          HollowCmd    = ProcessCommandLine,\n          HollowSHA256 = SHA256\n| order by Stage1Time desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet WindowMin = 5m;\nlet OfficeApps    = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"eqnedt32.exe\"]);\nlet ScriptHosts   = dynamic([\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"mshta.exe\"]);\nlet HollowTargets = dynamic([\"regasm.exe\",\"regsvcs.exe\",\"msbuild.exe\",\"installutil.exe\",\"aspnet_compiler.exe\",\"caspol.exe\"]);\nlet Stage1 = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | extend ParentFile = tolower(tostring(split(ParentProcessName, @\"\\\")[-1])),\n             ChildFile  = tolower(tostring(split(NewProcessName,    @\"\\\")[-1]))\n    | where ParentFile in (OfficeApps)\n    | where ChildFile  in (ScriptHosts)\n    | where CommandLine has_any (\"hidden\",\"-w 1\",\"-windowstyle\",\"FromBase64String\",\"DownloadString\",\"DownloadFile\",\"Invoke-Expression\",\" iex \",\"-EncodedCommand\",\"pastebin\",\"bit.ly\",\"j.mp\")\n    | project Stage1Time = TimeGenerated, Computer, Account,\n              OfficeParent = ParentFile, ScriptHost = ChildFile,\n              ScriptCmd    = CommandLine, ScriptPid = NewProcessId;\nSecurityEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventID == 4688\n| extend ChildFile  = tolower(tostring(split(NewProcessName,    @\"\\\")[-1])),\n         ParentFile = tolower(tostring(split(ParentProcessName, @\"\\\")[-1]))\n| where ChildFile  in (HollowTargets)\n| where ParentFile in (ScriptHosts)\n| join kind=inner Stage1 on Computer\n| where TimeGenerated between (Stage1Time .. Stage1Time + WindowMin)\n| project Computer, Account, Stage1Time, HollowTime = TimeGenerated,\n          DelaySec = datetime_diff('second', TimeGenerated, Stage1Time),\n          OfficeParent, ScriptHost, ScriptCmd,\n          HollowTarget = ChildFile, HollowCmd = CommandLine\n| order by Stage1Time desc",
      "sigma_yaml": "title: Gorgon Group .NET LOLBin Hollow-Target Spawned by Script Host\nid: 9b1c4d8e-7a2f-4b6e-9c11-7d3a4e5f6c20\nstatus: experimental\ndescription: Detects the second stage of the Gorgon Group infection chain \u2014 a .NET utility (RegAsm, RegSvcs, MSBuild, InstallUtil, aspnet_compiler) spawned by a Windows script host (wscript/cscript/powershell/mshta) used as a process-hollowing target for NJRat/Remcos/LokiBot/AZORult.\nreferences:\n  - https://attack.mitre.org/groups/G0078/\n  - https://unit42.paloaltonetworks.com/unit42-aggah-campaign-bit-ly-blogspot-and-pastebin-used-for-c2-in-large-scale-campaign/\nauthor: Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.execution\n  - attack.defense_evasion\n  - attack.t1055.012\n  - attack.t1059.001\n  - attack.t1564.003\n  - attack.g0078\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent:\n    ParentImage|endswith:\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\mshta.exe'\n  selection_child:\n    Image|endswith:\n      - '\\RegAsm.exe'\n      - '\\RegSvcs.exe'\n      - '\\MSBuild.exe'\n      - '\\InstallUtil.exe'\n      - '\\aspnet_compiler.exe'\n      - '\\caspol.exe'\n  condition: selection_parent and selection_child\nfalsepositives:\n  - Developer build automation invoking MSBuild from a script wrapper\n  - Custom in-house deployment tooling that scripts InstallUtil/RegAsm\nlevel: high",
      "rationale": "This isn't a generic 'Office spawns PowerShell' template \u2014 it specifically requires the THREE-stage Gorgon signature: (1) Office parent, (2) hidden-window script host with downloader / base64 / Bitly-Pastebin markers (Aggah campaign infrastructure), and (3) a .NET LOLBin spawned by that same script host within 5 minutes on the same device. RegAsm/MSBuild process hollowing is Gorgon's near-exclusiv"
    },
    {
      "title": "[LLM] Gorgon Group: Run-key persistence pointing to script-host payload in user-writable path",
      "description": "Hunts Gorgon Group's persistence step: after the initial PowerShell/VBScript stager executes, the actor writes an HKCU\\...\\Run or RunOnce value whose data invokes wscript/cscript/powershell/mshta against a payload staged in %APPDATA%\\Roaming, %TEMP%, %ProgramData% or %Public% \u2014 exactly the location set used by the Aggah/Subaat operators. Filtering on script-host + user-writable directory in the value data eliminates almost all legitimate vendor installers, which use signed EXE paths under Program Files.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "medium",
      "techs": [
        "T1547.001",
        "T1059.001",
        "T1059.005",
        "T1112"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true min(_time) as firstTime max(_time) as lastTime count from datamodel=Endpoint.Registry where (Registry.registry_path=\"*\\\\CurrentVersion\\\\Run\\\\*\" OR Registry.registry_path=\"*\\\\CurrentVersion\\\\RunOnce\\\\*\" OR Registry.registry_path=\"*\\\\Explorer\\\\Shell Folders\\\\Startup*\") AND Registry.action IN (\"created\",\"modified\") AND (Registry.registry_value_data=\"*wscript*\" OR Registry.registry_value_data=\"*cscript*\" OR Registry.registry_value_data=\"*powershell*\" OR Registry.registry_value_data=\"*pwsh*\" OR Registry.registry_value_data=\"*mshta*\") AND (Registry.registry_value_data=\"*\\\\AppData\\\\*\" OR Registry.registry_value_data=\"*\\\\Roaming\\\\*\" OR Registry.registry_value_data=\"*\\\\Temp\\\\*\" OR Registry.registry_value_data=\"*\\\\Public\\\\*\" OR Registry.registry_value_data=\"*\\\\ProgramData\\\\*\") by Registry.dest Registry.user Registry.registry_path Registry.registry_value_name Registry.registry_value_data Registry.process_name Registry.process_path | `drop_dm_object_name(Registry)` | where NOT (process_name IN (\"msiexec.exe\",\"explorer.exe\",\"setup.exe\",\"installer.exe\",\"squirrel.exe\",\"update.exe\")) | sort 0 - firstTime",
      "kql": "let LookbackDays = 14d;\nlet ScriptHosts = dynamic([\"wscript\",\"cscript\",\"powershell\",\"pwsh\",\"mshta\"]);\nlet UserWritable = dynamic([@\"\\AppData\\\", @\"\\Roaming\\\", @\"\\Local\\Temp\\\", @\"\\Public\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\"]);\nDeviceRegistryEvents\n| where Timestamp > ago(LookbackDays)\n| where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\", @\"\\Explorer\\Shell Folders\\Startup\", @\"\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\")\n| where RegistryValueData has_any (ScriptHosts)\n| where RegistryValueData has_any (UserWritable)\n| where InitiatingProcessAccountName !endswith \"$\"\n| where InitiatingProcessFileName !in~ (\"msiexec.exe\",\"setup.exe\",\"installer.exe\",\"squirrel.exe\",\"update.exe\")\n| extend SuspiciousScriptArg = iif(RegistryValueData has_any (\".vbs\",\".vbe\",\".js\",\".jse\",\".ps1\",\".hta\",\"-EncodedCommand\",\"FromBase64String\",\"-w 1\",\"hidden\"), true, false)\n| project Timestamp, DeviceName,\n          User = InitiatingProcessAccountName,\n          RegistryKey, RegistryValueName, RegistryValueData,\n          SuspiciousScriptArg,\n          WriterProcess = InitiatingProcessFileName,\n          WriterCmd     = InitiatingProcessCommandLine,\n          WriterParent  = InitiatingProcessParentFileName\n| order by Timestamp desc",
      "sentinel_kql": "let LookbackDays = 14d;\nlet ScriptHosts  = dynamic([\"wscript\",\"cscript\",\"powershell\",\"pwsh\",\"mshta\"]);\nlet UserWritable = dynamic([@\"\\AppData\\\", @\"\\Roaming\\\", @\"\\Local\\Temp\\\", @\"\\Public\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\"]);\nImRegistryEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\", @\"\\Explorer\\Shell Folders\\Startup\", @\"\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\")\n| where RegistryValueData has_any (ScriptHosts)\n| where RegistryValueData has_any (UserWritable)\n| where ActorUsername !endswith \"$\"\n| where ActingProcessName !in~ (\"msiexec.exe\",\"setup.exe\",\"installer.exe\",\"squirrel.exe\",\"update.exe\")\n| extend SuspiciousScriptArg = iif(RegistryValueData has_any (\".vbs\",\".vbe\",\".js\",\".jse\",\".ps1\",\".hta\",\"-EncodedCommand\",\"FromBase64String\",\"-w 1\",\"hidden\"), true, false)\n| project TimeGenerated, DvcHostname, ActorUsername, RegistryKey, RegistryValue, RegistryValueData,\n          SuspiciousScriptArg, ActingProcessName, ActingProcessCommandLine\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "Gorgon Group's persistence pattern is highly stylised \u2014 they almost never drop a signed EXE persistence; instead they chain HKCU\\Run \u2192 wscript/powershell \u2192 script under %APPDATA%\\Roaming or %ProgramData% (documented in Unit 42 'Aggah' and 'Subaat' write-ups, and replicated in 2019-2024 campaigns by Cybereason and Morphisec). Three filters make this fire only on the actor pattern: Run-family regist"
    }
  ]
}