# intel/.llm_actor_uc_cache/12/127afe90c07d1279c598724e323c71688cb6aca7.json
{
  "name": "Magic Hound",
  "ucs": [
    {
      "title": "[LLM] Magic Hound (APT35/Charming Kitten) CharmPower/PowerLess maldoc chain \u2014 Office-spawned scripting host followed by schtasks persistence",
      "description": "Hunts the canonical Magic Hound / TA453 / Phosphorus initial-access chain: a weaponised Office document (often a 'PDF/researcher invitation' lure aimed at journalists, academics, and policy staff) spawns PowerShell/WMIC/MSHTA which then registers a Scheduled Task to install the CharmPower or PowerLess .NET/PowerShell backdoor. The detection correlates an Office\u2192scripting-host descendant with a schtasks /create on the same device within a 10-minute window \u2014 a TTP combination repeatedly attributed to APT35 by Check Point, PwC, and Microsoft (CharmPower 2022, PowerLess 2022, BellaCiao 2023).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.001",
        "T1059.003",
        "T1047",
        "T1053.005",
        "T1027.010"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `summariesonly` min(_time) as stage_time from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\") AND Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"wmic.exe\",\"cmd.exe\") AND (Processes.process=\"*-enc*\" OR Processes.process=\"*EncodedCommand*\" OR Processes.process=\"*hidden*\" OR Processes.process=\"*FromBase64String*\" OR Processes.process=\"*DownloadString*\" OR Processes.process=\"*IEX*\" OR Processes.process=\"*Invoke-Expression*\" OR Processes.process=\"*Net.WebClient*\") by Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process Processes.process_guid | rename Processes.* as * | rename process_guid as stage_guid, process as stage_cmd, process_name as stage_image | join type=inner dest [ | tstats `summariesonly` min(_time) as persist_time values(Processes.process) as schedule_cmd from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" AND Processes.process=\"*/create*\" by Processes.dest Processes.user ] | where persist_time>=stage_time AND persist_time<=stage_time+600 | table stage_time persist_time dest user parent_process_name stage_image stage_cmd schedule_cmd",
      "kql": "// Magic Hound (APT35) CharmPower / PowerLess chain \u2014 Office macro -> scripting host -> schtasks persistence within 10 min\nlet win = 10m;\nlet officeApps = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"onenote.exe\"]);\nlet scriptHosts = dynamic([\"powershell.exe\",\"pwsh.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"wmic.exe\",\"cmd.exe\"]);\nlet stage1 =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where AccountName !endswith \"$\"\n    | where InitiatingProcessFileName in~ (officeApps)\n    | where FileName in~ (scriptHosts)\n    | where ProcessCommandLine has_any (\"-enc\",\"-Enc\",\"-EncodedCommand\",\"-w hidden\",\"-WindowStyle Hidden\",\"FromBase64String\",\"DownloadString\",\"DownloadFile\",\"IEX\",\"Invoke-Expression\",\"Net.WebClient\",\"WebRequest\",\"bitsadmin\")\n    | project StageTime = Timestamp, DeviceId, DeviceName, AccountName,\n              ParentImg = InitiatingProcessFileName,\n              StageImg  = FileName,\n              StageCmd  = ProcessCommandLine;\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has \"/create\"          // schtasks /create \u2014 Magic Hound persistence\n| where ProcessCommandLine has_any (\"/sc\",\"/SC\",\"/RU\",\"/TR\")\n| join kind=inner stage1 on DeviceId\n| where Timestamp between (StageTime .. StageTime + win)\n| project StageTime, PersistTime = Timestamp,\n          DelaySec = datetime_diff('second', Timestamp, StageTime),\n          DeviceName, AccountName,\n          ParentImg, StageImg, StageCmd,\n          ScheduleCmd = ProcessCommandLine\n| order by StageTime desc",
      "sentinel_kql": "// Magic Hound (APT35) CharmPower / PowerLess chain \u2014 Sentinel via SecurityEvent 4688\nlet win = 10m;\nlet officeApps = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"onenote.exe\"]);\nlet stage1 =\n    SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ParentBin = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1])),\n             ChildBin  = tolower(tostring(split(NewProcessName,    \"\\\\\")[-1]))\n    | where ParentBin in (officeApps)\n    | where ChildBin  in (\"powershell.exe\",\"pwsh.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"wmic.exe\",\"cmd.exe\")\n    | where CommandLine has_any (\"-enc\",\"-EncodedCommand\",\"-w hidden\",\"-WindowStyle Hidden\",\"FromBase64String\",\"DownloadString\",\"DownloadFile\",\"IEX\",\"Invoke-Expression\",\"Net.WebClient\")\n    | project StageTime = TimeGenerated, Computer, Account, ParentBin, ChildBin, StageCmd = CommandLine, StageLogonId = SubjectLogonId;\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where NewProcessName endswith \"\\\\schtasks.exe\"\n| where CommandLine has \"/create\"\n| where CommandLine has_any (\"/sc \",\"/SC \",\"/RU \",\"/TR \")\n| join kind=inner stage1 on Computer\n| where TimeGenerated between (StageTime .. StageTime + win)\n| project StageTime, PersistTime = TimeGenerated,\n          DelaySec = datetime_diff('second', TimeGenerated, StageTime),\n          Computer, Account, ParentBin, ChildBin, StageCmd, ScheduleCmd = CommandLine\n| order by StageTime desc",
      "sigma_yaml": "",
      "rationale": "Magic Hound's CharmPower (Check Point, Jan 2022), PowerLess (Cybereason, Feb 2022) and the BellaCiao implant (Bitdefender, Apr 2023) families share a near-identical bootstrapping chain: a benign-looking PDF/Word lure delivered through Newscaster-style social-engineering personas drops a stager that uses obfuscated PowerShell or WMIC out of the Office process tree, then immediately calls schtasks.e"
    },
    {
      "title": "[LLM] Magic Hound (TA453/Charming Kitten) post-compromise mailbox-rule abuse \u2014 auto-forward / delete inbox rules with external recipient",
      "description": "Hunts the credential-harvest \u2192 mailbox-takeover stage that defines Magic Hound operations against journalists, dissidents, and policy researchers (FBI/CISA AA22-264A, PwC 'Yellow Garuda', Mandiant). After phishing M365/Google credentials via spoofed login pages and bypassing MFA with token-relay kits, the actor logs in and creates Outlook inbox rules that silently forward correspondence to an attacker-controlled mailbox and/or auto-delete security warnings \u2014 a tradecraft Microsoft and Proofpoint have publicly attributed to APT35/TA453 since 2018.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1098.002",
        "T1114.002",
        "T1564.008",
        "T1078.004",
        "T1566.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "`o365_management_activity` Workload=Exchange Operation IN (\"New-InboxRule\",\"Set-InboxRule\",\"UpdateInboxRules\") ResultStatus=Succeeded | spath input=Parameters | rename Parameters{}.Name as p_name, Parameters{}.Value as p_value | eval kv=mvzip(p_name,p_value,\"=\") | mvexpand kv | rex field=kv \"^(?<rule_param>[^=]+)=(?<rule_value>.*)$\" | eval is_forward=if(match(rule_param,\"(?i)ForwardTo|ForwardAsAttachmentTo|RedirectTo\"),1,0), is_delete=if(match(rule_param,\"(?i)DeleteMessage|MoveToFolder\") AND match(rule_value,\"(?i)true|RSS|Junk|Deleted|Archive\"),1,0), is_external=if(match(rule_value,\"(?i)@\") AND NOT match(rule_value,\"(?i)@yourcorp\\.com|@yourcorp\\.onmicrosoft\\.com\"),1,0) | stats min(_time) as first_seen values(rule_param) as rule_params values(rule_value) as rule_values sum(is_forward) as fwd_hits sum(is_delete) as del_hits sum(is_external) as ext_hits values(ClientIP) as src_ip values(UserAgent) as ua by UserId Operation | where (fwd_hits>0 AND ext_hits>0) OR (fwd_hits>0 AND del_hits>0) | join type=left UserId [ | tstats `summariesonly` latest(_time) as last_signin values(Authentication.src) as signin_ip from datamodel=Authentication where Authentication.app=\"Office365\" Authentication.action=success by Authentication.user | rename Authentication.user as UserId ] | table first_seen UserId Operation rule_params rule_values src_ip ua signin_ip",
      "kql": "// Magic Hound (APT35/TA453) post-takeover inbox-rule abuse \u2014 auto-forward to external + hide-traces\nlet corp_domains = dynamic([\"yourcorp.com\",\"yourcorp.onmicrosoft.com\"]);   // <-- tune\nlet rule_ops = dynamic([\"New-InboxRule\",\"Set-InboxRule\",\"UpdateInboxRules\",\"Set-Mailbox\"]);\nCloudAppEvents\n| where Timestamp > ago(30d)\n| where Application in~ (\"Microsoft Exchange Online\",\"Office 365\")\n| where ActionType in~ (rule_ops)\n| extend Raw = parse_json(tostring(RawEventData))\n| extend Params = Raw.Parameters\n| mv-expand Params\n| extend ParamName  = tostring(Params.Name),\n         ParamValue = tostring(Params.Value)\n| where isnotempty(ParamName)\n| extend IsForwardParam = ParamName in~ (\"ForwardTo\",\"ForwardAsAttachmentTo\",\"RedirectTo\"),\n         IsHideParam    = (ParamName =~ \"DeleteMessage\" and ParamValue =~ \"True\") or\n                          (ParamName =~ \"MoveToFolder\"  and ParamValue has_any (\"RSS\",\"Junk\",\"Deleted\",\"Archive\",\"Conversation History\")),\n         IsKeywordRule  = ParamName in~ (\"SubjectContainsWords\",\"BodyContainsWords\",\"FromAddressContainsWords\")\n                          and ParamValue has_any (\"password\",\"login\",\"verification\",\"security\",\"sign-in\",\"sign in\",\"otp\",\"code\",\"alert\",\"phish\",\"helpdesk\")\n| extend ExternalRecipient = IsForwardParam\n                             and ParamValue matches regex @\"(?i)@\"\n                             and not (ParamValue has_any (corp_domains))\n| summarize FirstSeen = min(Timestamp),\n            RuleParams  = make_set(ParamName, 50),\n            RuleValues  = make_set(ParamValue, 50),\n            ForwardHits  = countif(IsForwardParam),\n            ExternalHits = countif(ExternalRecipient),\n            HideHits     = countif(IsHideParam),\n            KeywordHits  = countif(IsKeywordRule),\n            SrcIPs       = make_set(IPAddress, 10),\n            Countries    = make_set(CountryCode, 10),\n            UAs          = make_set(UserAgent, 10)\n            by AccountObjectId, AccountDisplayName, ActionType\n| where (ForwardHits > 0 and ExternalHits > 0)            // forward to outside org\n     or (ForwardHits > 0 and HideHits > 0)                // forward + hide-traces\n     or (HideHits    > 0 and KeywordHits > 0)             // delete-on-keyword (security-mail suppression)\n| join kind=leftouter (\n    AADSignInEventsBeta\n    | where Timestamp > ago(30d)\n    | where ErrorCode == 0\n    | summarize SignInIPs = make_set(IPAddress, 10),\n                SignInCountries = make_set(Country, 10),\n                RiskMax = max(tostring(RiskLevelDuringSignIn))\n                by AccountObjectId\n  ) on AccountObjectId\n| order by FirstSeen desc",
      "sentinel_kql": "// Magic Hound (APT35/TA453) post-takeover inbox-rule abuse \u2014 Sentinel via OfficeActivity\nlet corp_domains = dynamic([\"yourcorp.com\",\"yourcorp.onmicrosoft.com\"]);   // <-- tune\nOfficeActivity\n| where TimeGenerated > ago(30d)\n| where OfficeWorkload =~ \"Exchange\"\n| where Operation in~ (\"New-InboxRule\",\"Set-InboxRule\",\"UpdateInboxRules\",\"Set-Mailbox\")\n| where ResultStatus =~ \"Succeeded\" or isempty(ResultStatus)\n| extend P = parse_json(Parameters)\n| mv-expand P\n| extend ParamName  = tostring(P.Name),\n         ParamValue = tostring(P.Value)\n| extend IsForwardParam = ParamName in~ (\"ForwardTo\",\"ForwardAsAttachmentTo\",\"RedirectTo\"),\n         IsHideParam    = (ParamName =~ \"DeleteMessage\" and ParamValue =~ \"True\") or\n                          (ParamName =~ \"MoveToFolder\"  and ParamValue has_any (\"RSS\",\"Junk\",\"Deleted\",\"Archive\",\"Conversation History\")),\n         IsKeywordRule  = ParamName in~ (\"SubjectContainsWords\",\"BodyContainsWords\",\"FromAddressContainsWords\")\n                          and ParamValue has_any (\"password\",\"login\",\"verification\",\"security\",\"sign-in\",\"sign in\",\"otp\",\"code\",\"alert\",\"phish\",\"helpdesk\")\n| extend ExternalRecipient = IsForwardParam\n                             and ParamValue matches regex @\"(?i)@\"\n                             and not (ParamValue has_any (corp_domains))\n| summarize FirstSeen   = min(TimeGenerated),\n            RuleParams  = make_set(ParamName, 50),\n            RuleValues  = make_set(ParamValue, 50),\n            ForwardHits  = countif(IsForwardParam),\n            ExternalHits = countif(ExternalRecipient),\n            HideHits     = countif(IsHideParam),\n            KeywordHits  = countif(IsKeywordRule),\n            SrcIPs       = make_set(ClientIP, 10),\n            UAs          = make_set(UserAgent, 10)\n            by UserId, Operation\n| where (ForwardHits > 0 and ExternalHits > 0)\n     or (ForwardHits > 0 and HideHits     > 0)\n     or (HideHits    > 0 and KeywordHits  > 0)\n| join kind=leftouter (\n    SigninLogs\n    | where TimeGenerated > ago(30d)\n    | where ResultType == 0\n    | extend Country = tostring(parse_json(LocationDetails).countryOrRegion)\n    | summarize SignInIPs       = make_set(IPAddress, 10),\n                SignInCountries = make_set(Country, 10),\n                RiskMax         = max(tostring(RiskLevelDuringSignIn))\n                by UserPrincipalName\n  ) on $left.UserId == $right.UserPrincipalName\n| order by FirstSeen desc",
      "sigma_yaml": "",
      "rationale": "FBI/CISA AA22-264A, Microsoft DSU/MTAC ('Mint Sandstorm'/Phosphorus), Proofpoint TA453 reporting and Mandiant all describe the same post-credential-phish playbook for Magic Hound: log in from a residential/Iranian VPS, create a hidden Outlook rule that auto-forwards every inbound message to an attacker-controlled mailbox and silently moves security/MFA notifications to RSS/Deleted Items so the use"
    }
  ]
}