# rule_packs/sigma/UC_SCHEDULED_TASK.yml
title: Scheduled task created with suspicious image / encoded args
id: UC_SCHEDULED_TASK
status: experimental
description: schtasks.exe /create or Microsoft-Windows-TaskScheduler EventID 4698 with LOLBin actions.
references:
  - https://clankerusecase.com/use_cases/install/UC_SCHEDULED_TASK.yml
author: usecaseintel auto-generator
tags:
  - usecaseintel.tier.hunting
  - usecaseintel.kill_chain.install
  - attack.t1053_005
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|contains:
      - 'schtasks.exe'
      - '/create'
      - 'powershell'
      - 'cmd.exe'
      - 'rundll32'
      - '-enc'
      - 'FromBase64'
      - '\\Users\\Public'
      - '\\AppData\\'
  condition: selection
falsepositives:
  - Legitimate use of any of the above strings on dev / admin hosts
  - Tune via known-good allowlist before alerting
level: low
