# intel/.llm_actor_uc_cache/4b/4b7b1bb067fcb70a58f15d4ada892a7f265b0add.json
{
  "name": "RansomHub",
  "ucs": [
    {
      "title": "[LLM] RansomHub EDRKillShifter BYOVD: kernel service registers driver outside System32",
      "description": "Detects RansomHub affiliates dropping a vulnerable signed driver (EDRKillShifter / 'Loader.exe' / 'KillSecv1.exe') to a writable path and registering a kernel-mode service pointing at it \u2014 the BYOVD step before the loaded driver tears down EDR from kernel mode. This is RansomHub's signature pre-encryption defense-evasion move (also observed in attacks attributed to former ALPHV/BlackCat affiliates).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1014",
        "T1068",
        "T1543.003",
        "T1562.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=false count min(_time) as firstTime max(_time) as lastTime values(Registry.process_guid) as proc values(Registry.user) as user from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\SYSTEM\\\\CurrentControlSet\\\\Services\\\\*\" AND Registry.registry_value_name=\"ImagePath\" AND Registry.registry_value_data=\"*.sys\" AND NOT (Registry.registry_value_data IN (\"*\\\\System32\\\\drivers\\\\*\",\"*\\\\SystemRoot\\\\System32\\\\drivers\\\\*\",\"*\\\\Windows\\\\System32\\\\DriverStore\\\\*\")) AND (Registry.registry_value_data=\"*\\\\ProgramData\\\\*\" OR Registry.registry_value_data=\"*\\\\Users\\\\Public\\\\*\" OR Registry.registry_value_data=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Registry.registry_value_data=\"*\\\\Windows\\\\Temp\\\\*\" OR Registry.registry_value_data=\"*\\\\PerfLogs\\\\*\" OR Registry.registry_value_data=\"*\\\\$Recycle.Bin\\\\*\") by host Registry.registry_path Registry.registry_value_data | `drop_dm_object_name(Registry)` | rename host as dest | join type=left dest proc [| tstats summariesonly=false values(Processes.process_name) as loader_name values(Processes.process) as loader_cmd values(Processes.process_hash) as loader_hash from datamodel=Endpoint.Processes by host Processes.process_guid | `drop_dm_object_name(Processes)` | rename host as dest Processes.process_guid as proc] | convert ctime(firstTime) ctime(lastTime) | table firstTime dest user registry_path registry_value_data loader_name loader_cmd loader_hash",
      "kql": "let LookbackDays = 7d;\nlet _legit_driver_paths = dynamic([@\"\\system32\\drivers\\\", @\"\\systemroot\\system32\\drivers\\\", @\"\\windows\\system32\\driverstore\\\"]);\nlet _suspicious_paths = dynamic([@\"\\programdata\\\", @\"\\users\\public\\\", @\"\\appdata\\local\\temp\\\", @\"\\appdata\\roaming\\\", @\"\\windows\\temp\\\", @\"\\perflogs\\\", @\"\\$recycle.bin\\\"]);\nDeviceRegistryEvents\n| where Timestamp > ago(LookbackDays)\n| where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has @\"\\SYSTEM\\CurrentControlSet\\Services\\\"\n| where RegistryValueName =~ \"ImagePath\"\n| extend ImagePathLower = tolower(tostring(RegistryValueData))\n| where ImagePathLower endswith \".sys\"\n| where not(ImagePathLower has_any (_legit_driver_paths))\n| where ImagePathLower has_any (_suspicious_paths)\n| where InitiatingProcessAccountName !endswith \"$\"\n| where InitiatingProcessFileName !in~ (\"msiexec.exe\",\"trustedinstaller.exe\",\"wusa.exe\",\"setup.exe\")\n| extend NewServiceName = tostring(extract(@\"\\\\Services\\\\([^\\\\]+)\", 1, RegistryKey))\n| project Timestamp, DeviceName,\n          Loader = InitiatingProcessFileName,\n          LoaderPath = InitiatingProcessFolderPath,\n          LoaderCmd = InitiatingProcessCommandLine,\n          LoaderSHA256 = InitiatingProcessSHA256,\n          ActorAccount = InitiatingProcessAccountName,\n          NewServiceName,\n          DriverImagePath = RegistryValueData\n| order by Timestamp desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet _legit_driver_paths = dynamic([@\"\\system32\\drivers\\\", @\"\\systemroot\\system32\\drivers\\\", @\"\\windows\\system32\\driverstore\\\"]);\nlet _suspicious_paths = dynamic([@\"\\programdata\\\", @\"\\users\\public\\\", @\"\\appdata\\local\\temp\\\", @\"\\windows\\temp\\\", @\"\\perflogs\\\", @\"\\$recycle.bin\\\"]);\nEvent\n| where TimeGenerated > ago(LookbackDays)\n| where Source =~ \"Service Control Manager\" and EventID == 7045\n| extend SvcName    = extract(@\"(?i)Service Name:\\s*([^\\r\\n]+)\", 1, RenderedDescription),\n         ImagePath  = tolower(extract(@\"(?i)Service File Name:\\s*([^\\r\\n]+)\", 1, RenderedDescription)),\n         SvcType    = extract(@\"(?i)Service Type:\\s*([^\\r\\n]+)\", 1, RenderedDescription),\n         StartType  = extract(@\"(?i)Service Start Type:\\s*([^\\r\\n]+)\", 1, RenderedDescription),\n         SvcAccount = extract(@\"(?i)Service Account:\\s*([^\\r\\n]+)\", 1, RenderedDescription)\n| where (SvcType has \"kernel\") or ImagePath endswith \".sys\"\n| where not(ImagePath has_any (_legit_driver_paths))\n| where ImagePath has_any (_suspicious_paths)\n| project TimeGenerated, Computer, SvcName, ImagePath, SvcType, StartType, SvcAccount, UserName\n| order by TimeGenerated desc",
      "sigma_yaml": "title: BYOVD - Kernel Service Registers Driver Outside System32 (RansomHub EDRKillShifter)\nid: 8d7a3b1c-2f4e-49a6-bb12-1f5e0a7d9c3e\nstatus: experimental\ndescription: Detects creation of a kernel-mode service whose ImagePath points to a .sys driver outside the standard Windows driver store. Used by RansomHub's EDRKillShifter loader and other BYOVD tooling to terminate EDR from kernel mode.\nreferences:\n  - https://www.cisa.gov/news-events/cybersecurity-advisories/aa24-242a\n  - https://www.sophos.com/en-us/content/security-operations/active-adversary-report\nauthor: detection-engineering\ndate: 2026/05/13\ntags:\n  - attack.defense_evasion\n  - attack.privilege_escalation\n  - attack.t1014\n  - attack.t1068\n  - attack.t1543.003\n  - attack.t1562.001\nlogsource:\n  product: windows\n  category: registry_set\ndetection:\n  selection:\n    TargetObject|contains: '\\SYSTEM\\CurrentControlSet\\Services\\'\n    TargetObject|endswith: '\\ImagePath'\n    Details|endswith: '.sys'\n  suspect_path:\n    Details|contains:\n      - '\\ProgramData\\'\n      - '\\Users\\Public\\'\n      - '\\AppData\\Local\\Temp\\'\n      - '\\AppData\\Roaming\\'\n      - '\\Windows\\Temp\\'\n      - '\\PerfLogs\\'\n      - '\\$Recycle.Bin\\'\n  filter_legit:\n    Details|contains:\n      - '\\System32\\drivers\\'\n      - '\\SystemRoot\\System32\\drivers\\'\n      - '\\Windows\\System32\\DriverStore\\'\n  filter_installers:\n    Image|endswith:\n      - '\\msiexec.exe'\n      - '\\TrustedInstaller.exe'\n      - '\\wusa.exe'\n  condition: selection and suspect_path and not (filter_legit or filter_installers)\nfalsepositives:\n  - Niche third-party kernel drivers installed via non-MSI installers (rare; should be allowlisted by hash)\nlevel: high",
      "rationale": "EDRKillShifter is RansomHub's signature defense-evasion tool. Unlike generic AV/EDR-kill detections, this UC anchors on the structural step \u2014 kernel service registration whose ImagePath points outside the legit Windows driver store. Legit drivers are dropped to %SystemRoot%\\System32\\drivers\\ by MSI or TrustedInstaller; EDRKillShifter (and the same TTP used by Quantum, Akira, Medusa BYOVD variants "
    },
    {
      "title": "[LLM] RansomHub pre-encryption staging chain: VSS deletion + recovery disable + Defender tamper + service-stop fan-out within 10 minutes",
      "description": "Detects the tight, sequential pre-encryption playbook RansomHub affiliates execute on each victim host immediately before the Go-based encryptor runs: vssadmin/wmic shadow-copy deletion, bcdedit boot-recovery disable, Defender tamper (sc/PowerShell/reg), and a fan-out of net/sc/taskkill against backup, SQL, Veeam, EDR, and Defender services \u2014 all within a ~10-minute window from a single user account.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1490",
        "T1489",
        "T1562.001",
        "T1486",
        "T1070.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=false count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as commands from datamodel=Endpoint.Processes where Processes.user!=\"*$\" AND ((Processes.process_name=\"vssadmin.exe\" AND Processes.process=\"*delete*\" AND Processes.process=\"*shadow*\") OR (Processes.process_name=\"wmic.exe\" AND Processes.process=\"*shadowcopy*\" AND Processes.process=\"*delete*\") OR (Processes.process_name=\"bcdedit.exe\" AND (Processes.process=\"*recoveryenabled*No*\" OR Processes.process=\"*bootstatuspolicy*ignoreallfailures*\")) OR (Processes.process_name=\"sc.exe\" AND Processes.process=\"*config*\" AND Processes.process=\"*WinDefend*\" AND Processes.process=\"*disabled*\") OR (Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\") AND (Processes.process=\"*Set-MpPreference*DisableRealtimeMonitoring*\" OR Processes.process=\"*Add-MpPreference*ExclusionPath*\" OR Processes.process=\"*Win32_Shadowcopy*\")) OR (Processes.process_name=\"reg.exe\" AND Processes.process=\"*Windows Defender*\" AND Processes.process=\"*DisableAntiSpyware*\") OR (Processes.process_name IN (\"net.exe\",\"sc.exe\",\"taskkill.exe\") AND (Processes.process=\"*MSSQL*\" OR Processes.process=\"*SQLAgent*\" OR Processes.process=\"*Veeam*\" OR Processes.process=\"*Acronis*\" OR Processes.process=\"*BackupExec*\" OR Processes.process=\"*sophos*\" OR Processes.process=\"*mcafee*\" OR Processes.process=\"*CrowdStrike*\" OR Processes.process=\"*SentinelAgent*\" OR Processes.process=\"*MsMpSvc*\" OR Processes.process=\"*WinDefend*\"))) by host Processes.user _time span=10m | `drop_dm_object_name(Processes)` | eval is_vss=if(match(commands,\"(?i)(vssadmin.*delete.*shadow|wmic.*shadowcopy.*delete|Win32_Shadowcopy)\"),1,0), is_bcd=if(match(commands,\"(?i)bcdedit.*(recoveryenabled.*no|bootstatuspolicy.*ignoreallfailures)\"),1,0), is_def=if(match(commands,\"(?i)(WinDefend.*disabled|DisableRealtimeMonitoring|DisableAntiSpyware|ExclusionPath)\"),1,0), is_svcstop=if(match(commands,\"(?i)(net|sc|taskkill).*(MSSQL|Veeam|Acronis|BackupExec|sophos|mcafee|CrowdStrike|SentinelAgent|MsMpSvc|WinDefend)\"),1,0) | stats sum(is_vss) as vss sum(is_bcd) as bcd sum(is_def) as def sum(is_svcstop) as svcstop values(commands) as all_commands min(firstTime) as firstTime max(lastTime) as lastTime by host user | where vss>=1 AND ((bcd+def+svcstop)>=2) | convert ctime(firstTime) ctime(lastTime) | rename host as dest",
      "kql": "let LookbackHours = 24h;\nlet WindowMin = 10m;\nlet _backup_terms = dynamic([\"MSSQL\",\"SQLAgent\",\"SQLBrowser\",\"SQLWriter\",\"Veeam\",\"Backup\",\"Acronis\",\"BackupExec\",\"Postgres\",\"Oracle\",\"Exchange\",\"VSS\",\"BITS\",\"sophos\",\"mcafee\",\"symantec\",\"CrowdStrike\",\"SentinelAgent\",\"WinDefend\",\"MsMpSvc\"]);\nlet Vss = DeviceProcessEvents\n  | where Timestamp > ago(LookbackHours)\n  | where AccountName !endswith \"$\"\n  | where (FileName =~ \"vssadmin.exe\" and ProcessCommandLine has_all (\"delete\",\"shadow\"))\n       or (FileName =~ \"wmic.exe\" and ProcessCommandLine has_all (\"shadowcopy\",\"delete\"))\n       or (FileName in~ (\"powershell.exe\",\"pwsh.exe\") and ProcessCommandLine has \"Win32_Shadowcopy\")\n  | project VssTime = Timestamp, DeviceId, DeviceName, AccountName, VssCmd = ProcessCommandLine;\nlet Bcd = DeviceProcessEvents\n  | where Timestamp > ago(LookbackHours)\n  | where FileName =~ \"bcdedit.exe\"\n  | where ProcessCommandLine has_any (\"recoveryenabled no\",\"recoveryenabled No\",\"bootstatuspolicy ignoreallfailures\")\n  | project BcdTime = Timestamp, DeviceId, BcdCmd = ProcessCommandLine;\nlet Def = DeviceProcessEvents\n  | where Timestamp > ago(LookbackHours)\n  | where (FileName =~ \"sc.exe\" and ProcessCommandLine has_all (\"config\",\"WinDefend\",\"disabled\"))\n       or (FileName in~ (\"powershell.exe\",\"pwsh.exe\") and ProcessCommandLine has_any (\"Set-MpPreference -DisableRealtimeMonitoring\",\"Add-MpPreference -ExclusionPath\",\"Uninstall-WindowsFeature Windows-Defender\"))\n       or (FileName =~ \"reg.exe\" and ProcessCommandLine has \"Windows Defender\" and ProcessCommandLine has \"DisableAntiSpyware\")\n  | project DefTime = Timestamp, DeviceId, DefCmd = ProcessCommandLine;\nlet Stop = DeviceProcessEvents\n  | where Timestamp > ago(LookbackHours)\n  | where (FileName =~ \"net.exe\" and ProcessCommandLine has \"stop\")\n       or (FileName =~ \"sc.exe\" and ProcessCommandLine has \"stop\")\n       or (FileName =~ \"taskkill.exe\")\n  | where ProcessCommandLine has_any (_backup_terms)\n  | summarize StopCount = count(), StoppedSample = make_set(ProcessCommandLine, 10),\n              StopFirst = min(Timestamp), StopLast = max(Timestamp)\n              by DeviceId, DeviceName, bin(Timestamp, WindowMin)\n  | where StopCount >= 5;\nVss\n| join kind=inner Stop on DeviceId\n| where StopFirst between (VssTime - WindowMin .. VssTime + WindowMin)\n| join kind=leftouter Bcd on DeviceId\n| join kind=leftouter Def on DeviceId\n| where (isnotempty(BcdCmd) or isnotempty(DefCmd))\n| extend StagesObserved = strcat(\"VSS\",\n                              iif(isnotempty(BcdCmd),\",BCD\",\"\"),\n                              iif(isnotempty(DefCmd),\",DEF\",\"\"),\n                              \",SVC_STOP\")\n| project FirstActivity = VssTime, DeviceName, AccountName,\n          VssCmd, BcdCmd, DefCmd, StopCount, StoppedSample, StagesObserved\n| order by FirstActivity desc",
      "sentinel_kql": "let LookbackHours = 24h;\nlet WindowMin = 10m;\nlet _backup_terms = dynamic([\"MSSQL\",\"SQLAgent\",\"Veeam\",\"Backup\",\"Acronis\",\"BackupExec\",\"Exchange\",\"VSS\",\"sophos\",\"mcafee\",\"symantec\",\"CrowdStrike\",\"SentinelAgent\",\"WinDefend\",\"MsMpSvc\"]);\nlet Proc = SecurityEvent\n  | where TimeGenerated > ago(LookbackHours)\n  | where EventID == 4688\n  | where AccountType == \"User\" and SubjectUserName !endswith \"$\";\nlet Vss = Proc\n  | where (NewProcessName endswith \"\\\\vssadmin.exe\" and CommandLine has_all (\"delete\",\"shadow\"))\n       or (NewProcessName endswith \"\\\\wmic.exe\" and CommandLine has_all (\"shadowcopy\",\"delete\"))\n       or (NewProcessName has_any (\"\\\\powershell.exe\",\"\\\\pwsh.exe\") and CommandLine has \"Win32_Shadowcopy\")\n  | project VssTime = TimeGenerated, Computer, SubjectUserName, VssCmd = CommandLine;\nlet Bcd = Proc\n  | where NewProcessName endswith \"\\\\bcdedit.exe\"\n  | where CommandLine has_any (\"recoveryenabled no\",\"recoveryenabled No\",\"bootstatuspolicy ignoreallfailures\")\n  | project BcdTime = TimeGenerated, Computer, BcdCmd = CommandLine;\nlet Def = Proc\n  | where (NewProcessName endswith \"\\\\sc.exe\" and CommandLine has_all (\"config\",\"WinDefend\",\"disabled\"))\n       or (NewProcessName has_any (\"\\\\powershell.exe\",\"\\\\pwsh.exe\") and CommandLine has_any (\"Set-MpPreference -DisableRealtimeMonitoring\",\"Add-MpPreference -ExclusionPath\"))\n       or (NewProcessName endswith \"\\\\reg.exe\" and CommandLine has \"Windows Defender\" and CommandLine has \"DisableAntiSpyware\")\n  | project DefTime = TimeGenerated, Computer, DefCmd = CommandLine;\nlet Stop = Proc\n  | where (NewProcessName endswith \"\\\\net.exe\" and CommandLine has \"stop\")\n       or (NewProcessName endswith \"\\\\sc.exe\" and CommandLine has \"stop\")\n       or (NewProcessName endswith \"\\\\taskkill.exe\")\n  | where CommandLine has_any (_backup_terms)\n  | summarize StopCount = count(), StoppedSample = make_set(CommandLine, 10),\n              StopFirst = min(TimeGenerated)\n              by Computer, bin(TimeGenerated, WindowMin)\n  | where StopCount >= 5;\nVss\n| join kind=inner Stop on Computer\n| where StopFirst between (VssTime - WindowMin .. VssTime + WindowMin)\n| join kind=leftouter Bcd on Computer\n| join kind=leftouter Def on Computer\n| where (isnotempty(BcdCmd) or isnotempty(DefCmd))\n| project FirstActivity = VssTime, Computer, SubjectUserName,\n          VssCmd, BcdCmd, DefCmd, StopCount, StoppedSample\n| order by FirstActivity desc",
      "sigma_yaml": "",
      "rationale": "Individually, each leg of this chain (vssadmin delete shadows, bcdedit /set recoveryenabled No, sc config WinDefend disabled, net stop Veeam*) appears across many ransomware families \u2014 but the TIGHT 10-minute correlation of VSS-delete + (recovery-disable OR Defender-tamper) + 5+ backup/EDR service stops from a single user account on a single host is a near-fingerprint for RansomHub's pre-encryptio"
    }
  ]
}