# intel/.llm_actor_uc_cache/05/0582ebbf1021dcaeee191d3f320ef89d307efa86.json
{
  "name": "APT10",
  "ucs": [
    {
      "title": "[LLM] APT10 (Stone Panda / menuPass) DLL side-loading via signed Oracle/McAfee/VMware host binaries",
      "description": "Hunts APT10's signature initial-execution chain: drop a legitimately-signed third-party binary (jjs.exe, GUP.exe, mcoemcpy.exe, EhStorAuthn.exe, vmtoolsd.exe, etc.) into a non-install path together with a malicious sibling DLL (jli.dll, libcurl.dll, McUtil.dll, vmtools.dll) so the signed host loads APT10 implants like ChChes, RedLeaves, ANEL, UPPERCUT and PlugX. The host-binary + DLL pairings, together with the run-from-Temp/Public/ProgramData footprint, are documented across PwC, FireEye, JPCERT and Trend Micro reporting on Cloudhopper / Operation ChessMaster.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.002",
        "T1036.005",
        "T1027",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=false count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.process_path) as path from datamodel=Endpoint.Processes where Processes.process_name IN (\"jjs.exe\",\"GUP.exe\",\"mcoemcpy.exe\",\"EhStorAuthn.exe\",\"rasphone.exe\",\"chrome_frame_helper.exe\",\"vmtoolsd.exe\",\"mscoree.exe\",\"sllauncher.exe\") by host Processes.user Processes.process_name Processes.process_path Processes.process_id\n| `drop_dm_object_name(Processes)`\n| eval suspect_path=if(match(process_path,\"(?i)\\\\\\\\(Users|Public|ProgramData|AppData|Temp|Windows\\\\\\\\Temp|Intel|PerfLogs)\\\\\\\\\") AND NOT match(process_path,\"(?i)\\\\\\\\Program Files( \\\\(x86\\\\))?\\\\\\\\\"),1,0)\n| where suspect_path=1\n| join type=inner host process_id [\n    | tstats summariesonly=false count values(Filesystem.file_path) as dll_path values(Filesystem.file_hash) as dll_hash from datamodel=Endpoint.Filesystem where Filesystem.file_name IN (\"jli.dll\",\"libcurl.dll\",\"McUtil.dll\",\"vmtools.dll\",\"EhStorShell.dll\",\"chrome_frame_helper.dll\",\"mscoree.dll\",\"VERSION.dll\") by host Filesystem.process_id Filesystem.file_name\n    | `drop_dm_object_name(Filesystem)`\n    | rename process_id as process_id\n  ]\n| convert ctime(firstTime) ctime(lastTime)\n| table firstTime lastTime host user process_name process_path cmdline parent file_name dll_path dll_hash",
      "kql": "// APT10 DLL side-loading: signed host binary in a non-install path loads the matching APT10 sideload DLL adjacent to it\nlet SideloadHosts = dynamic([\"jjs.exe\",\"gup.exe\",\"mcoemcpy.exe\",\"ehstorauthn.exe\",\"rasphone.exe\",\"chrome_frame_helper.exe\",\"vmtoolsd.exe\",\"mscoree.exe\",\"sllauncher.exe\"]);\nlet SideloadDlls   = dynamic([\"jli.dll\",\"libcurl.dll\",\"mcutil.dll\",\"vmtools.dll\",\"ehstorshell.dll\",\"chrome_frame_helper.dll\",\"version.dll\"]);\nlet KnownGoodPaths = dynamic([@\"c:\\program files\\\", @\"c:\\program files (x86)\\\", @\"c:\\windows\\system32\\\", @\"c:\\windows\\syswow64\\\"]);\nDeviceImageLoadEvents\n| where Timestamp > ago(14d)\n| where InitiatingProcessFileName in~ (SideloadHosts)\n| where FileName in~ (SideloadDlls)\n// adjacent-DLL is the actual sideload pattern \u2014 DLL written next to the loader\n| extend HostDir = tolower(InitiatingProcessFolderPath), DllDir = tolower(FolderPath)\n| where HostDir == DllDir\n// strip legitimate vendor install paths\n| where not(HostDir startswith_cs @\"c:\\program files\\\")\n     and not(HostDir startswith_cs @\"c:\\program files (x86)\\\")\n     and not(HostDir startswith_cs @\"c:\\windows\\\")\n| project Timestamp, DeviceName,\n          User = InitiatingProcessAccountName,\n          HostBin = InitiatingProcessFileName,\n          HostBinPath = InitiatingProcessFolderPath,\n          HostBinSha256 = InitiatingProcessSHA256,\n          HostBinCmdLine = InitiatingProcessCommandLine,\n          LoadedDll = FileName,\n          LoadedDllPath = FolderPath,\n          LoadedDllSha256 = SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// APT10 DLL side-loading via Sysmon EID 7 (ImageLoad). Requires Sysmon \u2192 WindowsEvent\nlet SideloadHosts = dynamic([\"jjs.exe\",\"gup.exe\",\"mcoemcpy.exe\",\"ehstorauthn.exe\",\"rasphone.exe\",\"chrome_frame_helper.exe\",\"vmtoolsd.exe\",\"mscoree.exe\",\"sllauncher.exe\"]);\nlet SideloadDlls   = dynamic([\"jli.dll\",\"libcurl.dll\",\"mcutil.dll\",\"vmtools.dll\",\"ehstorshell.dll\",\"chrome_frame_helper.dll\",\"version.dll\"]);\nWindowsEvent\n| where TimeGenerated > ago(14d)\n| where Provider == \"Microsoft-Windows-Sysmon\" and EventID == 7\n| extend Image       = tolower(tostring(EventData.Image)),\n         ImageLoaded = tolower(tostring(EventData.ImageLoaded)),\n         User        = tostring(EventData.User),\n         SignerName  = tostring(EventData.SignatureStatus)\n| extend HostBin  = extract(@\"\\\\([^\\\\]+)$\", 1, Image),\n         LoadedDll = extract(@\"\\\\([^\\\\]+)$\", 1, ImageLoaded),\n         HostDir   = extract(@\"^(.+)\\\\[^\\\\]+$\", 1, Image),\n         DllDir    = extract(@\"^(.+)\\\\[^\\\\]+$\", 1, ImageLoaded)\n| where HostBin in (SideloadHosts) and LoadedDll in (SideloadDlls)\n| where HostDir == DllDir\n| where not(HostDir startswith @\"c:\\program files\\\")\n     and not(HostDir startswith @\"c:\\program files (x86)\\\")\n     and not(HostDir startswith @\"c:\\windows\\\")\n| project TimeGenerated, Computer, User, HostBin, Image, LoadedDll, ImageLoaded\n| order by TimeGenerated desc",
      "sigma_yaml": "title: APT10 (Stone Panda) Signed-Binary DLL Side-Loading Pair\nid: 9e3f4a18-1c4f-4f2c-8c1c-0b3a4d27a01b\nstatus: experimental\ndescription: Detects APT10 / menuPass DLL side-loading where a known signed third-party host binary loads the matching malicious DLL pair from a non-install path. Pairings are taken from PwC, FireEye, JPCERT and Trend Micro reporting on RedLeaves, ANEL, UPPERCUT, ChChes and PlugX.\nreferences:\n  - https://www.pwc.co.uk/issues/cyber-security-services/insights/operation-cloud-hopper.html\n  - https://www.jpcert.or.jp/english/at/2017/at170024.html\nauthor: Detection Engineering\ndate: 2026/05/05\nlogsource:\n  category: image_load\n  product: windows\ndetection:\n  selection_pair:\n    Image|endswith:\n      - '\\jjs.exe'\n      - '\\GUP.exe'\n      - '\\mcoemcpy.exe'\n      - '\\EhStorAuthn.exe'\n      - '\\rasphone.exe'\n      - '\\chrome_frame_helper.exe'\n      - '\\vmtoolsd.exe'\n      - '\\mscoree.exe'\n      - '\\sllauncher.exe'\n    ImageLoaded|endswith:\n      - '\\jli.dll'\n      - '\\libcurl.dll'\n      - '\\McUtil.dll'\n      - '\\vmtools.dll'\n      - '\\EhStorShell.dll'\n      - '\\chrome_frame_helper.dll'\n      - '\\version.dll'\n  filter_legit_install:\n    Image|contains:\n      - '\\Program Files\\'\n      - '\\Program Files (x86)\\'\n      - '\\Windows\\System32\\'\n      - '\\Windows\\SysWOW64\\'\n  condition: selection_pair and not filter_legit_install\nfalsepositives:\n  - Vendor utilities installed to a non-default folder (rare in managed estates)\nlevel: high\ntags:\n  - attack.defense_evasion\n  - attack.t1574.002\n  - attack.t1036.005\n  - attack.g0045",
      "rationale": "Generic 'unsigned DLL load by lsass' or 'rundll32 from Temp' rules will miss APT10 \u2014 they specifically pair a signed Oracle/McAfee/VMware/Trend host binary with a sibling DLL that hijacks an exported function. The (host-binary, DLL-name, adjacent-folder) tuple is what makes this APT10-specific rather than a generic side-loading template; the same pairings have been re-used across the Cloudhopper, "
    },
    {
      "title": "[LLM] APT10 Cloudhopper post-compromise: AD recon (csvde/net/dsquery) chained to remote scheduled-task creation on managed hosts",
      "description": "Hunts the APT10 Cloudhopper hand-off pattern: an MSP/jump-host operator account runs csvde / net group / dsquery / nltest to enumerate the customer domain, then within ~30 minutes uses schtasks /create /s <remoteHost> (or wmic /node:) to drop a scheduled task on a downstream client host \u2014 the documented APT10 mechanism for pivoting from MSP infrastructure to managed-customer endpoints. Tying the recon and the remote-schtasks call together on the same actor and host eliminates the high-volume noise of either signal alone.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1087.002",
        "T1018",
        "T1053.005",
        "T1021.002",
        "T1199"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=false count min(_time) as reconTime from datamodel=Endpoint.Processes where (Processes.process_name=\"csvde.exe\" OR Processes.process_name=\"dsquery.exe\" OR Processes.process_name=\"nltest.exe\" OR (Processes.process_name IN (\"net.exe\",\"net1.exe\") AND (Processes.process=\"*group \\\"domain admins\\\"*\" OR Processes.process=\"*group \\\"enterprise admins\\\"*\" OR Processes.process=\"*/domain*\"))) AND NOT Processes.user IN (\"*$\",\"SYSTEM\",\"LOCAL SERVICE\",\"NETWORK SERVICE\") by host Processes.user Processes.process_name Processes.process\n| `drop_dm_object_name(Processes)`\n| rename process as recon_cmd, process_name as recon_bin\n| join type=inner host user [\n    | tstats summariesonly=false count min(_time) as schtasksTime values(Processes.process) as schtasks_cmd from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" AND Processes.process=\"*/create*\" AND Processes.process=\"*/s *\" by host Processes.user\n    | `drop_dm_object_name(Processes)`\n  ]\n| eval delta=abs(schtasksTime-reconTime)\n| where delta<=1800\n| eval remote_host=replace(schtasks_cmd,\"(?i).*?/s\\s+\\\\\\\\?\\\\\\\\?([\\w\\.-]+).*\",\"\\1\")\n| where remote_host!=host\n| convert ctime(reconTime) ctime(schtasksTime)\n| table reconTime schtasksTime host user recon_bin recon_cmd schtasks_cmd remote_host delta",
      "kql": "// APT10 Cloudhopper recon-then-remote-schtasks chain on the same actor / pivot host\nlet LookbackDays = 14d;\nlet WindowMin    = 30;\nlet Recon =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where (FileName =~ \"csvde.exe\")\n         or (FileName in~ (\"net.exe\",\"net1.exe\")\n             and ProcessCommandLine has_any (\"group \\\"domain admins\\\"\",\"group \\\"enterprise admins\\\"\",\"/domain\"))\n         or (FileName =~ \"dsquery.exe\" and ProcessCommandLine has_any (\"user\",\"computer\",\"group\"))\n         or (FileName =~ \"nltest.exe\" and ProcessCommandLine has \"/dclist\")\n         or (FileName =~ \"ntdsutil.exe\")\n    | project ReconTime = Timestamp, DeviceName, AccountName,\n              ReconBin = FileName, ReconCmd = ProcessCommandLine;\nDeviceProcessEvents\n| where Timestamp > ago(LookbackDays)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has \"/create\" and ProcessCommandLine has \"/s \"\n| where AccountName !endswith \"$\"\n| extend RemoteHost = tolower(extract(@\"(?i)/s\\s+\\\\\\\\?([\\w\\.-]+)\", 1, ProcessCommandLine))\n| where isnotempty(RemoteHost) and RemoteHost != tolower(DeviceName)\n| join kind=inner Recon on DeviceName, AccountName\n| where ReconTime between (Timestamp - WindowMin * 1m .. Timestamp + WindowMin * 1m)\n| project SchtasksTime = Timestamp, ReconTime,\n          DelaySec = datetime_diff('second', Timestamp, ReconTime),\n          DeviceName, AccountName, IsInitiatingProcessRemoteSession,\n          ReconBin, ReconCmd,\n          SchtasksCmd = ProcessCommandLine,\n          RemoteHost,\n          Parent = InitiatingProcessFileName,\n          ParentCmd = InitiatingProcessCommandLine\n| order by SchtasksTime desc",
      "sentinel_kql": "// APT10 Cloudhopper recon \u2192 remote schtasks chain via ASIM ImProcessCreate\nlet LookbackDays = 14d;\nlet WindowMin    = 30;\nlet Recon =\n    ImProcessCreate\n    | where TimeGenerated > ago(LookbackDays)\n    | where ActorUsername !endswith \"$\"\n    | where (TargetProcessName has \"csvde.exe\")\n         or (TargetProcessName has_any (\"net.exe\",\"net1.exe\")\n             and TargetProcessCommandLine has_any (\"group \\\"domain admins\\\"\",\"group \\\"enterprise admins\\\"\",\"/domain\"))\n         or (TargetProcessName has \"dsquery.exe\" and TargetProcessCommandLine has_any (\"user\",\"computer\",\"group\"))\n         or (TargetProcessName has \"nltest.exe\" and TargetProcessCommandLine has \"/dclist\")\n         or (TargetProcessName has \"ntdsutil.exe\")\n    | project ReconTime = TimeGenerated, DvcHostname, ActorUsername,\n              ReconBin = TargetProcessName, ReconCmd = TargetProcessCommandLine;\nImProcessCreate\n| where TimeGenerated > ago(LookbackDays)\n| where TargetProcessName has \"schtasks.exe\"\n| where TargetProcessCommandLine has \"/create\" and TargetProcessCommandLine has \"/s \"\n| where ActorUsername !endswith \"$\"\n| extend RemoteHost = tolower(extract(@\"(?i)/s\\s+\\\\\\\\?([\\w\\.-]+)\", 1, TargetProcessCommandLine))\n| where isnotempty(RemoteHost) and RemoteHost != tolower(DvcHostname)\n| join kind=inner Recon on DvcHostname, ActorUsername\n| where ReconTime between (TimeGenerated - WindowMin * 1m .. TimeGenerated + WindowMin * 1m)\n| project SchtasksTime = TimeGenerated, ReconTime,\n          DelaySec = datetime_diff('second', TimeGenerated, ReconTime),\n          DvcHostname, ActorUsername, ReconBin, ReconCmd,\n          SchtasksCmd = TargetProcessCommandLine, RemoteHost,\n          Parent = ParentProcessName, ParentCmd = ParentProcessCommandLine\n| order by SchtasksTime desc",
      "sigma_yaml": "",
      "rationale": "The Cloudhopper / Operation TradeSecret reporting (PwC, NCSC-UK, US-CERT TA17-117A) repeatedly attribute the same hand-off shape to APT10: csvde-style domain enumeration from an MSP jump host, immediately followed by remote scheduled-task creation against managed-customer endpoints over RPC/SMB. Either signal alone is dashboard-tier \u2014 IT runs schtasks /s and admins occasionally csvde \u2014 but the sam",
      "_field_issues": [
        {
          "kind": "unknown_field",
          "field": "IsInitiatingProcessRemoteSession",
          "tables_in_scope": [
            "DeviceProcessEvents"
          ],
          "suggestion": "InitiatingProcessTokenElevation",
          "message": "`IsInitiatingProcessRemoteSession` is not a canonical column of DeviceProcessEvents. Did you mean `InitiatingProcessTokenElevation`?"
        }
      ]
    }
  ]
}