# intel/.llm_actor_uc_cache/43/434ea2866b63c422db188d3832b79be3f85def2a.json
{
  "name": "Rancor",
  "ucs": [
    {
      "title": "[LLM] Rancor (G0075) DDKONG/PLAINTEE installer: Office/WScript spawning msiexec to fetch remote MSI",
      "description": "Hunts Rancor's signature delivery chain where a politically-themed lure document (typically targeting ASEAN/SE-Asia government and diplomatic entities) drops a VBScript that invokes msiexec /q /i against a remote HTTP(S) URL to stage DDKONG or PLAINTEE. The Office/wscript-->msiexec-->external-URL combination is a near-fingerprint of Rancor's 2018+ campaigns documented by Unit 42.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.005",
        "T1059.003",
        "T1218.007",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process) as parent_cmd values(Processes.user) as user values(Processes.process_id) as pid from datamodel=Endpoint.Processes where Processes.process_name=msiexec.exe (Processes.parent_process_name IN (winword.exe,excel.exe,powerpnt.exe,outlook.exe,wscript.exe,cscript.exe,mshta.exe,cmd.exe)) (Processes.process=\"*http://*\" OR Processes.process=\"*https://*\" OR Processes.process=\"*\\\\\\\\*\") (Processes.process=\"* /q*\" OR Processes.process=\"* /quiet*\" OR Processes.process=\"* /i *\") by Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process_guid | `drop_dm_object_name(Processes)` | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` | join type=left process_guid [| tstats summariesonly=true values(All_Traffic.dest) as remote_dest values(All_Traffic.dest_port) as dest_port from datamodel=Network_Traffic.All_Traffic where All_Traffic.app=msiexec.exe by All_Traffic.process_guid | rename All_Traffic.process_guid as process_guid] | where isnotnull(remote_dest)",
      "kql": "let suspiciousParents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"cmd.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(30d)\n| where FileName =~ \"msiexec.exe\"\n| where InitiatingProcessFileName has_any (suspiciousParents)\n| where ProcessCommandLine has_any (\"http://\",\"https://\",\"\\\\\\\\\")\n| where ProcessCommandLine has_any (\"/q\",\"/quiet\",\"/i \")\n| extend Chain = strcat(InitiatingProcessParentFileName, \" -> \", InitiatingProcessFileName, \" -> \", FileName)\n| join kind=leftouter (\n    DeviceNetworkEvents\n    | where Timestamp > ago(30d)\n    | where InitiatingProcessFileName =~ \"msiexec.exe\"\n    | where RemoteUrl has_any (\".asia\",\".ph\",\".kh\",\".sg\",\".my\",\".id\",\".vn\",\".th\") or RemoteIPType == \"Public\"\n    | summarize RemoteUrls=make_set(RemoteUrl,10), RemoteIPs=make_set(RemoteIP,10) by DeviceId, InitiatingProcessId, InitiatingProcessCreationTime\n) on DeviceId, $left.ProcessId == $right.InitiatingProcessId\n| project Timestamp, DeviceName, AccountName, Chain, ProcessCommandLine, InitiatingProcessCommandLine, RemoteUrls, RemoteIPs, ReportId\n| order by Timestamp desc",
      "rationale": "Rancor's Unit 42-documented chain (the Cambodia/Singapore lures) is uniquely characterised by Office or Windows Script Host invoking msiexec with /q /i and a remote HTTP path to stage the DDKONG loader and PLAINTEE backdoor. Legitimate enterprise software almost never uses msiexec to silently fetch an MSI from an external URL via a script host or Office child, making the parent-child + URL combina"
    },
    {
      "title": "[LLM] Rancor PLAINTEE post-install persistence: schtasks + Reg Run + WMI subscription writing C2 over HTTP",
      "description": "Hunts Rancor's PLAINTEE post-compromise persistence pattern, where the implant registers itself via schtasks /create or a Run-key while beaconing to hard-coded C2 over HTTP/UDP. The combination of scheduled-task or AccessibilityFeatures-style debugger persistence created by a non-standard binary in %APPDATA% / %TEMP%, followed by outbound HTTP from that same binary, mirrors the PLAINTEE behaviour observed in Rancor's SE-Asia operations.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "medium",
      "techs": [
        "T1053.005",
        "T1546.003",
        "T1071.001",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.process_path) as path from datamodel=Endpoint.Processes where (Processes.process_name=schtasks.exe AND Processes.process=\"*/create*\" AND (Processes.process=\"*\\\\AppData\\\\*\" OR Processes.process=\"*\\\\Temp\\\\*\" OR Processes.process=\"*\\\\ProgramData\\\\*\")) OR (Processes.process_name=wmic.exe AND Processes.process=\"*__EventFilter*\") by Processes.dest Processes.user Processes.process_guid Processes.process_name | `drop_dm_object_name(Processes)` | join type=inner dest [| tstats summariesonly=true values(All_Traffic.dest) as c2_dest values(All_Traffic.dest_port) as c2_port values(All_Traffic.app) as c2_app from datamodel=Network_Traffic.All_Traffic where (All_Traffic.app=\"*\\\\AppData\\\\*\" OR All_Traffic.app=\"*\\\\Temp\\\\*\" OR All_Traffic.app=\"*\\\\ProgramData\\\\*\") (All_Traffic.dest_port IN (80,443,8080,53)) NOT All_Traffic.dest_category=internal by All_Traffic.src | rename All_Traffic.src as dest] | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "let userDirs = dynamic([\"\\\\AppData\\\\\",\"\\\\Temp\\\\\",\"\\\\ProgramData\\\\\",\"\\\\Public\\\\\"]);\nlet persist =\n  DeviceProcessEvents\n  | where Timestamp > ago(14d)\n  | where (FileName =~ \"schtasks.exe\" and ProcessCommandLine has \"/create\" and ProcessCommandLine has_any (userDirs))\n     or (FileName =~ \"wmic.exe\" and ProcessCommandLine has \"__EventFilter\")\n     or (FileName =~ \"reg.exe\" and ProcessCommandLine has_any (\"\\\\CurrentVersion\\\\Run\",\"Image File Execution Options\") and ProcessCommandLine has_any (userDirs))\n  | project DeviceId, DeviceName, persistTime=Timestamp, persistCmd=ProcessCommandLine, persistParent=InitiatingProcessFileName, persistImage=InitiatingProcessFolderPath;\nlet beacon =\n  DeviceNetworkEvents\n  | where Timestamp > ago(14d)\n  | where RemoteIPType == \"Public\"\n  | where RemotePort in (80, 443, 8080, 53)\n  | where InitiatingProcessFolderPath has_any (userDirs)\n  | where InitiatingProcessFileName !in~ (\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"OneDrive.exe\",\"Teams.exe\")\n  | summarize beaconCount=count(), beaconHosts=dcount(RemoteIP), firstBeacon=min(Timestamp), lastBeacon=max(Timestamp), c2=make_set(RemoteUrl, 5), ips=make_set(RemoteIP, 5) by DeviceId, InitiatingProcessFileName, InitiatingProcessFolderPath\n  | where beaconCount > 5;\npersist\n| join kind=inner beacon on DeviceId\n| project DeviceName, persistTime, persistCmd, InitiatingProcessFileName, InitiatingProcessFolderPath, beaconCount, c2, ips, firstBeacon, lastBeacon\n| order by persistTime desc",
      "rationale": "Pairing persistence creation (schtasks/WMI subscription/Run-key) sourced from a user-writable directory with subsequent web-protocol beaconing from the same staged binary models the documented PLAINTEE implant lifecycle Rancor uses against SE-Asian government targets, rather than alerting on either persistence or C2 alone."
    }
  ]
}