# intel/.llm_actor_uc_cache/9b/9b0f88d18592277bc68e75397ede6d1c686c1085.json
{
  "name": "Putter Panda",
  "ucs": [
    {
      "title": "[LLM] Putter Panda (PLA 61486) MSUpdater/4H RAT persistence via Adobe-masquerading Run keys",
      "description": "Hunts the characteristic Putter Panda dropper chain documented by CrowdStrike: a stage-1 dropper decodes an encrypted payload to %APPDATA%/%LOCALAPPDATA% under filenames mimicking Adobe ARM / MSUpdater, then registers HKCU\\...\\Run persistence pointing at that decoded binary. Combines T1027.013 (encoded payload on disk) with T1547.001 (Run key) in one chain.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1027.013",
        "T1547.001",
        "T1036.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Registry.registry_value_data) as run_target values(Registry.process_name) as writer values(Registry.user) as user from datamodel=Endpoint.Registry where Registry.registry_path IN (\"*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run*\",\"*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce*\") AND (Registry.registry_value_name IN (\"AdobeUpdater\",\"AdobeARM\",\"Adobe ARM\",\"MSUpdate\",\"MSUpdater\",\"AcroRd32\") OR Registry.registry_value_data IN (\"*\\\\AppData\\\\*\\\\AdobeARM*\",\"*\\\\AppData\\\\*\\\\MSUpdate*\",\"*\\\\ProgramData\\\\Adobe\\\\ARM\\\\*\",\"*\\\\AppData\\\\Roaming\\\\Adobe\\\\Update*\")) by Registry.dest Registry.user Registry.registry_value_name Registry.registry_value_data \n| `drop_dm_object_name(Registry)` \n| join type=inner dest [| tstats summariesonly=t count from datamodel=Endpoint.Filesystem where (Filesystem.file_path IN (\"*\\\\AppData\\\\*\\\\AdobeARM*\",\"*\\\\AppData\\\\*\\\\MSUpdate*\",\"*\\\\AppData\\\\Roaming\\\\Adobe\\\\Update*\")) AND Filesystem.action=\"created\" by Filesystem.dest Filesystem.file_path Filesystem.process_name | `drop_dm_object_name(Filesystem)` | rename process_name as dropper] \n| where dropper!=\"AdobeARM.exe\" AND dropper!=\"armsvc.exe\" AND NOT match(dropper,\"(?i)reader_sl|acrord32|setup\") \n| eval risk_score=85, mitre=\"T1547.001+T1027.013\", actor=\"Putter Panda / APT2 (G0024)\" \n| table firstTime lastTime dest user writer dropper run_target file_path actor mitre risk_score",
      "kql": "let suspectPaths = pack_array(@\"\\AppData\\Roaming\\Adobe\\Update\\\", @\"\\AppData\\Local\\Adobe\\ARM\\\", @\"\\AppData\\Roaming\\MSUpdate\\\", @\"\\AppData\\Local\\MSUpdater\\\", @\"\\ProgramData\\Adobe\\ARM\\\");\nlet drops = DeviceFileEvents\n| where Timestamp > ago(7d)\n| where ActionType == \"FileCreated\"\n| where FolderPath has_any (suspectPaths)\n      or (FileName matches regex @\"(?i)^(AdobeARM|MSUpdate|MSUpdater|AcroRd32)\\.(exe|dll)$\" and FolderPath startswith \"C:\\\\Users\\\\\")\n| where InitiatingProcessFolderPath !startswith @\"C:\\Program Files\\Adobe\\\"\n      and InitiatingProcessFolderPath !startswith @\"C:\\Program Files (x86)\\Adobe\\\"\n      and InitiatingProcessFolderPath !startswith @\"C:\\Program Files (x86)\\Common Files\\Adobe\\\"\n| project DropTime=Timestamp, DeviceId, DeviceName, DroppedPath=FolderPath, DroppedFile=FileName, DropperProcess=InitiatingProcessFileName, DropperPath=InitiatingProcessFolderPath, DropperSha=InitiatingProcessSHA256, AccountName=InitiatingProcessAccountName;\nlet runkeys = DeviceRegistryEvents\n| where Timestamp > ago(7d)\n| where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has_any (@\"\\Microsoft\\Windows\\CurrentVersion\\Run\", @\"\\Microsoft\\Windows\\CurrentVersion\\RunOnce\")\n| where RegistryValueName in~ (\"AdobeUpdater\",\"AdobeARM\",\"Adobe ARM\",\"MSUpdate\",\"MSUpdater\",\"AcroRd32\")\n      or RegistryValueData has_any (@\"\\AppData\\Roaming\\Adobe\\Update\\\", @\"\\AppData\\Local\\Adobe\\ARM\\\", @\"\\AppData\\Roaming\\MSUpdate\\\", @\"\\AppData\\Local\\MSUpdater\\\")\n| where InitiatingProcessFolderPath !startswith @\"C:\\Program Files\\Adobe\\\"\n      and InitiatingProcessFolderPath !startswith @\"C:\\Program Files (x86)\\Adobe\\\"\n| project RunKeyTime=Timestamp, DeviceId, RegistryKey, RegistryValueName, RegistryValueData, RunKeyWriter=InitiatingProcessFileName;\ndrops\n| join kind=inner (runkeys) on DeviceId\n| where datetime_diff('minute', RunKeyTime, DropTime) between (-30 .. 60)\n| extend Actor=\"Putter Panda / APT2 (G0024) \u2014 PLA Unit 61486\", MitreChain=\"T1027.013 -> T1547.001\"\n| project DropTime, RunKeyTime, DeviceName, AccountName, DropperProcess, DropperPath, DropperSha, DroppedPath, DroppedFile, RegistryKey, RegistryValueName, RegistryValueData, Actor, MitreChain",
      "rationale": "Putter Panda's MSUpdater family (4H RAT, 3PARA RAT droppers) is documented by CrowdStrike's 2014 'Putter Panda' report as consistently masquerading as Adobe Reader Update Manager (AdobeARM) or Microsoft Update components, dropping decoded payloads to user-writable Adobe-themed paths and pinning persistence with HKCU Run values that legitimate Adobe/Microsoft installers never write. The temporal jo"
    },
    {
      "title": "[LLM] Putter Panda 3PARA RAT DLL injection into long-running browser/host processes",
      "description": "Hunts 3PARA RAT's signature execution pattern: a small stage-1 PE injects a malicious DLL (often named to mimic Microsoft components like msi.dll, winhlp.dll, ssdpapi.dll) into a long-running process such as iexplore.exe, svchost.exe, or explorer.exe to acquire C2 over HTTP. Maps T1055.001 (DLL Injection) \u2014 Putter Panda's primary execution-stealth move per CrowdStrike's reverse-engineering of 3PARA/4H.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "medium",
      "techs": [
        "T1055.001",
        "T1574.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.process_name IN (\"rundll32.exe\",\"regsvr32.exe\") AND (Processes.process IN (\"*\\\\AppData\\\\*\\\\msi.dll*\",\"*\\\\AppData\\\\*\\\\winhlp.dll*\",\"*\\\\AppData\\\\*\\\\ssdpapi.dll*\",\"*\\\\ProgramData\\\\*\\\\msi.dll*\",\"*\\\\Temp\\\\*\\\\msi.dll*\") OR match(Processes.process,\"(?i)rundll32.*,\\s*(Initialize|StartW|ServiceMain|InstallService)\")) by Processes.dest Processes.user Processes.process_name Processes.process Processes.parent_process_name \n| `drop_dm_object_name(Processes)` \n| join type=inner dest [| tstats summariesonly=t count values(All_Traffic.dest_ip) as c2_ip values(All_Traffic.dest_port) as c2_port from datamodel=Network_Traffic.All_Traffic where (All_Traffic.app=\"iexplore.exe\" OR All_Traffic.app=\"svchost.exe\" OR All_Traffic.app=\"rundll32.exe\") AND All_Traffic.dest_port IN (80,443,8080) AND NOT All_Traffic.dest_category=\"internal\" by All_Traffic.src All_Traffic.app | rename All_Traffic.src as dest] \n| eval risk_score=80, mitre=\"T1055.001+T1574.002\", actor=\"Putter Panda / APT2 (G0024) \u2014 3PARA/4H RAT\" \n| table firstTime lastTime dest user parent process_name cmdline c2_ip c2_port actor mitre risk_score",
      "kql": "let suspectDllNames = dynamic([\"msi.dll\",\"winhlp.dll\",\"ssdpapi.dll\",\"netapi.dll\",\"wininet.dll\",\"version.dll\"]);\nlet sideload = DeviceImageLoadEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (\"iexplore.exe\",\"svchost.exe\",\"explorer.exe\",\"rundll32.exe\",\"regsvr32.exe\")\n| where FileName in~ (suspectDllNames)\n| where FolderPath !startswith @\"C:\\Windows\\System32\\\"\n      and FolderPath !startswith @\"C:\\Windows\\SysWOW64\\\"\n      and FolderPath !startswith @\"C:\\Program Files\\\"\n| where FolderPath startswith @\"C:\\Users\\\" or FolderPath startswith @\"C:\\ProgramData\\\" or FolderPath has @\"\\Temp\\\"\n| project LoadTime=Timestamp, DeviceId, DeviceName, HostProcess=InitiatingProcessFileName, HostPath=InitiatingProcessFolderPath, MaliciousDll=FileName, DllPath=FolderPath, DllSha=SHA256, AccountName=InitiatingProcessAccountName;\nlet beacon = DeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (\"iexplore.exe\",\"svchost.exe\",\"rundll32.exe\")\n| where RemotePort in (80,443,8080)\n| where RemoteIPType == \"Public\"\n| summarize Beacons=count(), DistinctRemotes=dcount(RemoteIP), RemoteIPs=make_set(RemoteIP,10), RemoteUrls=make_set(RemoteUrl,10) by DeviceId, InitiatingProcessFileName, bin(Timestamp,1h)\n| where Beacons >= 5;\nsideload\n| join kind=inner (beacon) on DeviceId\n| where InitiatingProcessFileName == HostProcess\n| extend Actor=\"Putter Panda / APT2 (G0024) \u2014 3PARA RAT\", MitreChain=\"T1574.002 -> T1055.001 -> C2\"\n| project LoadTime, DeviceName, AccountName, HostProcess, HostPath, MaliciousDll, DllPath, DllSha, RemoteIPs, RemoteUrls, Beacons, Actor, MitreChain",
      "rationale": "CrowdStrike's reverse-engineering of 3PARA RAT specifically identified a small DLL (commonly named to collide with system DLLs like msi.dll) that is loaded by a host process from a non-system path and then beacons over HTTP. The pivot \u2014 a system-named DLL loading from %AppData%/%Temp%/%ProgramData% INTO a network-facing host process that then beacons \u2014 is exactly the Putter Panda execution chain a"
    }
  ]
}