# intel/.llm_actor_uc_cache/4e/4e190815b03aa0f15e4832f051fb8ab141a38844.json
{
  "name": "Darkhotel",
  "ucs": [
    {
      "title": "[LLM] Darkhotel (DUBNIUM) hotel-Wi-Fi update hijack \u2014 masqueraded signed 'updater' executed from user-writable path with Run-key persistence",
      "description": "Hunts Darkhotel's signature initial-access pattern: travelling executives are MITM'd on hotel/conference networks and served a trojanised software update (GoogleUpdate, AdobeARM, FlashPlayerUpdate, HangulUpdate, etc.) signed with a stolen or factored weak-key certificate, then dropped to %APPDATA%/%TEMP% with a Run-key persistence and Authenticode publisher mismatch. This combination of masquerade + non-vendor signer + temp execution + Run key is highly specific to Darkhotel and not seen in benign update flows.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1189",
        "T1036.005",
        "T1553.002",
        "T1547.001",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.parent_process_name) as parent values(Processes.process) as cmdline values(Processes.process_hash) as sha256 values(Processes.process_signature) as signer values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.process_name IN (\"GoogleUpdate.exe\",\"AdobeARM.exe\",\"AdobeUpdater.exe\",\"FlashPlayerUpdate*.exe\",\"HncUpdate*.exe\",\"swupdate*.exe\",\"jucheck.exe\",\"OfficeC2RClient.exe\") (Processes.process_path=\"*\\\\AppData\\\\*\" OR Processes.process_path=\"*\\\\Temp\\\\*\" OR Processes.process_path=\"*\\\\Downloads\\\\*\" OR Processes.process_path=\"*\\\\ProgramData\\\\*\") by host Processes.process_name Processes.process_path Processes.dest \n| `drop_dm_object_name(Processes)` \n| where NOT match(signer,\"(?i)Google LLC|Adobe Inc|Adobe Systems|Microsoft Corporation|Oracle America|Hancom|Sun Microsystems\") \n| join type=inner host [| tstats summariesonly=t values(Registry.registry_value_data) as runVal values(Registry.registry_path) as runPath from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\CurrentVersion\\\\Run*\" (Registry.registry_value_data=\"*\\\\AppData\\\\*\" OR Registry.registry_value_data=\"*\\\\Temp\\\\*\" OR Registry.registry_value_data=\"*\\\\ProgramData\\\\*\") by host] \n| join type=inner host [| tstats summariesonly=t values(All_Traffic.dest_ip) as priorDest values(All_Traffic.dest_port) as priorPort from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest_port IN (80,443,8080) earliest=-30m@m by host | rename All_Traffic.* as *] \n| eval signer_anomaly=if(isnull(signer) OR signer=\"\",\"unsigned\",\"non_vendor_signer\") \n| table firstTime host user process_name process_path sha256 signer signer_anomaly parent cmdline runPath runVal priorDest \n| `security_content_ctime(firstTime)`",
      "kql": "let updaterNames = dynamic([\"GoogleUpdate.exe\",\"AdobeARM.exe\",\"AdobeUpdater.exe\",\"FlashPlayerUpdate.exe\",\"FlashPlayerUpdateService.exe\",\"HncUpdate.exe\",\"swupdate.exe\",\"jucheck.exe\",\"OfficeC2RClient.exe\"]);\nlet vendorSigners = dynamic([\"Google LLC\",\"Google Inc\",\"Adobe Inc.\",\"Adobe Systems\",\"Microsoft Corporation\",\"Oracle America, Inc.\",\"Hancom Inc.\"]);\nlet susProcs = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName in~ (updaterNames)\n| where FolderPath has_any (@\"\\AppData\\\", @\"\\Temp\\\", @\"\\Downloads\\\", @\"\\ProgramData\\\")\n| extend BadSigner = iff(isempty(Signer) or Signer !in (vendorSigners), true, false)\n| where BadSigner == true\n| project Timestamp, DeviceId, DeviceName, AccountName, FileName, FolderPath, SHA256, Signer, SignatureState=tostring(parse_json(AdditionalFields).SignatureState), ProcessId=ProcessId, InitiatingProcessFileName, ProcessCommandLine;\nlet runKeyPersist = DeviceRegistryEvents\n| where Timestamp > ago(7d)\n| where ActionType == \"RegistryValueSet\"\n| where RegistryKey has @\"\\CurrentVersion\\Run\"\n| where RegistryValueData has_any (@\"\\AppData\\\", @\"\\Temp\\\", @\"\\ProgramData\\\")\n| project DeviceId, RunTime=Timestamp, RegistryKey, RegistryValueName, RegistryValueData;\nlet preNet = DeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where RemotePort in (80, 443, 8080)\n| where RemoteIPType !in (\"Private\",\"Loopback\",\"LinkLocal\")\n| summarize NetTime=min(Timestamp), DistinctRemoteIPs=dcount(RemoteIP), Sample=make_set(RemoteIP, 5) by DeviceId, InitiatingProcessFileName;\nsusProcs\n| join kind=inner runKeyPersist on DeviceId\n| where RunTime between ((Timestamp - 10m) .. (Timestamp + 30m))\n| join kind=leftouter preNet on DeviceId\n| where NetTime between ((Timestamp - 30m) .. Timestamp)\n| project Timestamp, DeviceName, AccountName, FileName, FolderPath, Signer, SHA256, ProcessCommandLine, RegistryKey, RegistryValueData, NetTime, Sample, DistinctRemoteIPs",
      "rationale": "Darkhotel's defining tradecraft (per Kaspersky's 2014 'Darkhotel APT' report and Microsoft's DUBNIUM tracking) is hijacking hotel networks to push fake updates signed with stolen/factored 512-bit RSA certificates from Korean CAs \u2014 i.e. the binary LOOKS legitimate (correct name, valid Authenticode chain) but the signer is NOT the vendor whose updater it is impersonating. Stacking (a) updater filena"
    },
    {
      "title": "[LLM] Darkhotel weaponised HWP/Office spearphish \u2192 cmstp/wscript dropper with anti-VM checks \u2192 encrypted C2 to long-dormant infra",
      "description": "Hunts Darkhotel's East-Asia-focused spearphishing chain (frequently Hangul HWP or Office docs targeting Korean/Japanese/Chinese executives) where macro/exploit execution spawns a script-host child that performs sandbox/VM evasion (process enumeration for vmtoolsd/vboxservice, system info collection, sleep/timing checks) and then beacons over TLS \u2014 the T1566.001 \u2192 T1203 \u2192 T1497 \u2192 T1573.001 chain Kaspersky and ESET have repeatedly attributed to Darkhotel/Tapaoux.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "exploit",
      "conf": "medium",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1203",
        "T1059.003",
        "T1497.001",
        "T1057",
        "T1082",
        "T1140"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime values(Processes.process) as cmdline values(Processes.process_path) as path values(Processes.process_hash) as sha256 from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"WINWORD.EXE\",\"EXCEL.EXE\",\"POWERPNT.EXE\",\"Hwp.exe\",\"Hword.exe\",\"HShell.exe\",\"HwpApp.exe\") Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"cmstp.exe\",\"rundll32.exe\",\"regsvr32.exe\") by host user Processes.parent_process_name Processes.process_name Processes.dest \n| `drop_dm_object_name(Processes)` \n| join type=inner host [| tstats summariesonly=t count as discoveryHits values(Processes.process) as discoveryCmd from datamodel=Endpoint.Processes where Processes.process_name IN (\"tasklist.exe\",\"wmic.exe\",\"systeminfo.exe\",\"reg.exe\",\"sc.exe\") (Processes.process=\"*vmtoolsd*\" OR Processes.process=\"*vboxservice*\" OR Processes.process=\"*vmware*\" OR Processes.process=\"*VirtualBox*\" OR Processes.process=\"*Sandboxie*\" OR Processes.process=\"*HKLM\\\\SYSTEM\\\\ControlSet*\\\\Services\\\\Disk*\" OR Processes.process=\"*BIOS*\" OR Processes.process=\"*ComputerSystemProduct*\") by host | where discoveryHits >=1] \n| join type=inner host [| tstats summariesonly=t values(All_Traffic.dest) as c2_dest values(All_Traffic.dest_port) as c2_port from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest_port IN (443,8443,4443) All_Traffic.app=\"ssl\" by host | rename All_Traffic.* as *] \n| join type=outer src_user [| tstats summariesonly=t values(All_Email.src_user) as recipient values(All_Email.file_name) as attach values(All_Email.subject) as subj from datamodel=Email.All_Email where All_Email.file_name IN (\"*.hwp\",\"*.hwpx\",\"*.doc\",\"*.docx\",\"*.xls\",\"*.xlsx\",\"*.lnk\",\"*.iso\",\"*.zip\") by All_Email.src_user | rename All_Email.src_user as src_user] \n| eval risk = 60 + if(parent_process_name LIKE \"%Hwp%\",20,0) + if(match(c2_port,\"4443|8443\"),10,0) \n| where risk >= 70 \n| table firstTime host user parent_process_name process_name path sha256 cmdline discoveryCmd c2_dest c2_port subj attach \n| `security_content_ctime(firstTime)`",
      "kql": "let officeParents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"hwp.exe\",\"hword.exe\",\"hshell.exe\",\"hwpapp.exe\"]);\nlet scriptHosts = dynamic([\"cmd.exe\",\"powershell.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"cmstp.exe\",\"rundll32.exe\",\"regsvr32.exe\"]);\nlet vmStrings = dynamic([\"vmtoolsd\",\"vboxservice\",\"vmware\",\"virtualbox\",\"sandboxie\",\"qemu\",\"xen\",\"ComputerSystemProduct\",\"VBoxGuest\"]);\nlet stage1 = DeviceProcessEvents\n| where Timestamp > ago(14d)\n| where InitiatingProcessFileName in~ (officeParents)\n| where FileName in~ (scriptHosts)\n| project S1Time=Timestamp, DeviceId, DeviceName, AccountName, ParentProc=InitiatingProcessFileName, FileName, ProcessCommandLine, S1Pid=ProcessId, S1Sha=SHA256;\nlet stage2 = DeviceProcessEvents\n| where Timestamp > ago(14d)\n| where FileName in~ (\"tasklist.exe\",\"wmic.exe\",\"systeminfo.exe\",\"reg.exe\",\"sc.exe\",\"powershell.exe\")\n| where ProcessCommandLine has_any (vmStrings) or ProcessCommandLine has_any (\"BIOS\",\"Manufacturer\",\"Model\",\"NumberOfCores\")\n| project S2Time=Timestamp, DeviceId, EvasionCmd=ProcessCommandLine, EvasionProc=FileName;\nlet stage3 = DeviceNetworkEvents\n| where Timestamp > ago(14d)\n| where RemotePort in (443, 8443, 4443)\n| where RemoteIPType !in (\"Private\",\"Loopback\")\n| project S3Time=Timestamp, DeviceId, RemoteIP, RemoteUrl, RemotePort, NetProc=InitiatingProcessFileName;\nlet phish = EmailAttachmentInfo\n| where Timestamp > ago(14d)\n| where FileName endswith \".hwp\" or FileName endswith \".hwpx\" or FileName endswith \".doc\" or FileName endswith \".docx\" or FileName endswith \".xls\" or FileName endswith \".xlsx\" or FileName endswith \".lnk\" or FileName endswith \".iso\"\n| join kind=inner (EmailEvents | where Timestamp > ago(14d)) on NetworkMessageId\n| project NetworkMessageId, RecipientEmailAddress, SenderFromAddress, SenderIPv4, Subject, FileName, FileType, ThreatTypes;\nstage1\n| join kind=inner stage2 on DeviceId\n| where S2Time between ((S1Time) .. (S1Time + 15m))\n| join kind=inner stage3 on DeviceId\n| where S3Time between ((S1Time) .. (S1Time + 30m))\n| extend UserUPN = tostring(AccountName)\n| join kind=leftouter (phish | extend UserUPN = tostring(split(RecipientEmailAddress,\"@\")[0])) on UserUPN\n| project S1Time, DeviceName, AccountName, ParentProc, FileName, ProcessCommandLine, EvasionProc, EvasionCmd, NetProc, RemoteIP, RemoteUrl, RemotePort, SenderFromAddress, Subject, FileName1=FileName1, FileType, ThreatTypes\n| sort by S1Time desc",
      "rationale": "Darkhotel/Tapaoux/DUBNIUM operations documented by Kaspersky (2014, 2015, 2017), ESET, and Microsoft repeatedly chain (1) HWP or Office spearphish targeting Korean/Japanese/Chinese executives, (2) child script-host or LOLBin (cmstp/wscript/regsvr32) carrying obfuscated payload, (3) VM/sandbox checks via WMI/registry/process enum (T1497.001 \u2014 Darkhotel's 'Tapaoux' and 'Asruex' samples are notorious"
    }
  ]
}