# intel/.llm_actor_uc_cache/4e/4e190815b03aa0f15e4832f051fb8ab141a38844.json
{
  "name": "Darkhotel",
  "ucs": [
    {
      "title": "[LLM] Darkhotel (DUBNIUM/Zigzag Hail) fake-software-update dropper from hotel/captive-portal staging path",
      "description": "Hunts the Darkhotel signature pattern of a digitally-signed binary masquerading as a legitimate updater (Flash, GoogleToolbar, Adobe, Realtek) that lands in a user-writable path after a hotel/Wi-Fi drive-by (T1189) and either sleeps to defeat sandboxes (T1497) or pulls a second-stage via PowerShell/cmd (T1059). Darkhotel famously hijacks hotel networks to inject these masqueraded updates against traveling executives, then signs payloads with stolen/forged certificates (T1553.002).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1189",
        "T1036.005",
        "T1553.002",
        "T1204.002",
        "T1497.003"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.process_path) as path values(Processes.process_hash) as sha256 from datamodel=Endpoint.Processes where Processes.user!=\"*SYSTEM\" Processes.user!=\"*$\" (Processes.process_name IN (\"*update*.exe\",\"*setup*.exe\",\"*flash*.exe\",\"*adobe*.exe\",\"*googletoolbar*.exe\",\"*hwp*.exe\",\"*realtek*.exe\")) (Processes.process_path=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Processes.process_path=\"*\\\\AppData\\\\Roaming\\\\*\" OR Processes.process_path=\"*\\\\Downloads\\\\*\" OR Processes.process_path=\"*\\\\Public\\\\*\" OR Processes.process_path=\"*\\\\ProgramData\\\\*\") by host Processes.user Processes.process_name Processes.process_path Processes.process_hash Processes.parent_process_name | rename Processes.* as * | join type=outer host process_hash [| tstats summariesonly=true count from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"*update*.exe\",\"*setup*.exe\",\"*flash*.exe\",\"*adobe*.exe\",\"*googletoolbar*.exe\",\"*hwp*.exe\",\"*realtek*.exe\") (Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"bitsadmin.exe\",\"certutil.exe\")) by host Processes.parent_process_name Processes.process_name Processes.process | rename Processes.parent_process_name as process_name Processes.process_name as child_process Processes.process as child_cmd] | where isnotnull(child_process) | convert ctime(firstTime) ctime(lastTime) | table firstTime host user process_name path process_hash parent child_process child_cmd",
      "kql": "let LookbackDays = 7d;\nlet UpdaterNames = dynamic([\"update\",\"setup\",\"flash\",\"adobe\",\"googletoolbar\",\"hwp\",\"realtek\",\"java\",\"silverlight\"]);\nlet UserPaths = dynamic([@\"\\AppData\\Local\\Temp\\\", @\"\\AppData\\Roaming\\\", @\"\\Downloads\\\", @\"\\Public\\\", @\"\\ProgramData\\\"]);\nlet LolBins = dynamic([\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"curl.exe\"]);\nlet Droppers = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where FolderPath has_any (UserPaths)\n    | where FileName has_any (UpdaterNames)\n    | where FileName endswith \".exe\"\n    // Darkhotel relocates signed binaries OR drops new ones whose version-info signer doesn't match Microsoft\n    | where InitiatingProcessVersionInfoCompanyName !has \"Microsoft\"\n    | project DropperTime = Timestamp, DeviceId, DeviceName, AccountName,\n              DropperFile = FileName, DropperPath = FolderPath,\n              DropperSha256 = SHA256, DropperPid = ProcessId,\n              DropperCompany = ProcessVersionInfoCompanyName,\n              DropperOriginalName = ProcessVersionInfoOriginalFileName,\n              DropperCmd = ProcessCommandLine;\nDeviceProcessEvents\n| where Timestamp > ago(LookbackDays)\n| where FileName in~ (LolBins)\n| where ProcessCommandLine has_any (\"http\", \"DownloadString\", \"IEX\", \"-enc\", \"-EncodedCommand\", \"FromBase64String\", \"WebClient\", \"BitsTransfer\", \"-w hidden\", \"-windowstyle hidden\")\n   or InitiatingProcessFileName has_any (UpdaterNames)\n| join kind=inner Droppers on $left.InitiatingProcessId == $right.DropperPid, DeviceId\n| where Timestamp between (DropperTime .. DropperTime + 10m)\n| project DropperTime, ChildTime = Timestamp, DeviceName, AccountName,\n          DropperFile, DropperPath, DropperSha256, DropperCompany, DropperOriginalName,\n          ChildProcess = FileName, ChildCmd = ProcessCommandLine,\n          DelaySec = datetime_diff('second', Timestamp, DropperTime)\n| order by DropperTime desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet UpdaterPattern = @\"(?i)\\\\(update|setup|flash|adobe|googletoolbar|hwp|realtek|java|silverlight)[^\\\\]*\\.exe$\";\nlet UserPathPattern = @\"(?i)\\\\(AppData\\\\Local\\\\Temp|AppData\\\\Roaming|Downloads|Public|ProgramData)\\\\\";\nlet LolBins = dynamic([\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"curl.exe\"]);\nlet Droppers = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | where NewProcessName matches regex UserPathPattern\n    | where NewProcessName matches regex UpdaterPattern\n    | project DropperTime = TimeGenerated, Computer, Account,\n              DropperImage = NewProcessName, DropperPid = NewProcessId,\n              DropperCmd = CommandLine;\nSecurityEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventID == 4688\n| extend ChildName = tostring(split(NewProcessName, \"\\\\\")[-1])\n| where tolower(ChildName) in (LolBins)\n| where CommandLine has_any (\"http\", \"DownloadString\", \"IEX\", \"-enc\", \"-EncodedCommand\", \"FromBase64String\", \"WebClient\", \"BitsTransfer\", \"-w hidden\", \"-windowstyle hidden\")\n| join kind=inner Droppers on Computer, $left.ProcessId == $right.DropperPid\n| where TimeGenerated between (DropperTime .. DropperTime + 10m)\n| project DropperTime, ChildTime = TimeGenerated, Computer, Account,\n          DropperImage, DropperCmd,\n          ChildImage = NewProcessName, ChildCmd = CommandLine,\n          DelaySec = datetime_diff('second', TimeGenerated, DropperTime)\n| order by DropperTime desc",
      "sigma_yaml": "",
      "rationale": "Darkhotel's distinctive delivery is a fake software-update binary planted via hotel Wi-Fi/captive-portal injection, and the dropper is almost always named after legit updaters (Flash/Adobe/Realtek/GoogleToolbar/HWP) \u2014 a pattern documented by Kaspersky and Microsoft. Generic 'update.exe in Temp' alone is noisy, but the actor-specific combination \u2014 masqueraded updater name + user-writable path + lol"
    },
    {
      "title": "[LLM] Darkhotel spearphish: Office/HWP decoy spawning binary masquerading as Windows system file with HKCU Run-key persistence",
      "description": "Catches Darkhotel's spearphishing chain (T1566.001 \u2192 T1203/T1204.002 \u2192 T1547.001) where a decoy document \u2014 frequently a Korean Hancom HWP, RTF, or weaponised Office file \u2014 exploits/spawns a payload that masquerades as a legitimate Windows binary (svchost.exe, igfxtray.exe, conhost.exe, wuauclt.exe, ctfmon.exe) but executes from a user profile directory, then writes a HKCU Run / RunOnce value within minutes for persistence. Darkhotel is documented to disguise loaders with these exact filenames and rely on user-writable persistence to survive without admin rights.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1203",
        "T1036.005",
        "T1547.001",
        "T1140"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as spawnTime values(Processes.process) as childCmd values(Processes.process_path) as childPath values(Processes.process_hash) as childHash values(Processes.process_id) as childPid from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"hwp.exe\",\"hword.exe\",\"acrord32.exe\",\"acrobat.exe\",\"eqnedt32.exe\") Processes.user!=\"*$\" Processes.process_name IN (\"svchost.exe\",\"igfxtray.exe\",\"conhost.exe\",\"wuauclt.exe\",\"ctfmon.exe\",\"taskhost.exe\",\"taskhostw.exe\",\"explorer.exe\",\"lsass.exe\",\"services.exe\",\"winlogon.exe\",\"smss.exe\") (Processes.process_path!=\"C:\\\\Windows\\\\System32\\\\*\" Processes.process_path!=\"C:\\\\Windows\\\\SysWOW64\\\\*\" Processes.process_path!=\"C:\\\\Windows\\\\WinSxS\\\\*\" Processes.process_path!=\"C:\\\\Program Files*\") by host Processes.user Processes.parent_process_name Processes.process_name Processes.process_path Processes.process_hash Processes.process_id | rename Processes.* as * | join type=inner host [| tstats summariesonly=true count min(_time) as regTime values(Registry.registry_key_name) as regKey values(Registry.registry_value_name) as regValueName values(Registry.registry_value_data) as regValueData values(Registry.process_id) as regPid from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\CurrentVersion\\\\Run*\" Registry.action=\"modified\" by host Registry.process_id | rename Registry.* as *] | where regPid==process_id OR (regTime>=spawnTime AND regTime<=spawnTime+600) | convert ctime(spawnTime) ctime(regTime) | table spawnTime regTime host user parent_process_name process_name process_path process_hash regKey regValueName regValueData",
      "kql": "let LookbackDays = 7d;\nlet PersistenceWindow = 15m;\nlet OfficeDecoyParents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"hwp.exe\",\"hword.exe\",\"hcell.exe\",\"hshow.exe\",\"acrord32.exe\",\"acrobat.exe\",\"eqnedt32.exe\"]);\nlet MasqueradeNames = dynamic([\"svchost.exe\",\"igfxtray.exe\",\"conhost.exe\",\"wuauclt.exe\",\"ctfmon.exe\",\"taskhost.exe\",\"taskhostw.exe\",\"explorer.exe\",\"lsass.exe\",\"services.exe\",\"winlogon.exe\",\"smss.exe\",\"spoolsv.exe\"]);\nlet MasqueradedSpawns = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where InitiatingProcessFileName in~ (OfficeDecoyParents)\n    | where FileName in~ (MasqueradeNames)\n    // genuine Windows binaries live under System32/SysWOW64/WinSxS \u2014 anywhere else is masquerading\n    | where not(FolderPath startswith @\"C:\\Windows\\System32\\\")\n    | where not(FolderPath startswith @\"C:\\Windows\\SysWOW64\\\")\n    | where not(FolderPath startswith @\"C:\\Windows\\WinSxS\\\")\n    | where not(FolderPath startswith @\"C:\\Program Files\\\")\n    | where not(FolderPath startswith @\"C:\\Program Files (x86)\\\")\n    | project SpawnTime = Timestamp, DeviceId, DeviceName, AccountName,\n              ParentDoc = InitiatingProcessFileName,\n              ParentDocCmd = InitiatingProcessCommandLine,\n              MasqFile = FileName, MasqPath = FolderPath,\n              MasqSha256 = SHA256, MasqPid = ProcessId,\n              MasqOriginalName = ProcessVersionInfoOriginalFileName,\n              MasqCompany = ProcessVersionInfoCompanyName;\nDeviceRegistryEvents\n| where Timestamp > ago(LookbackDays)\n| where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\", @\"\\Explorer\\Run\", @\"\\Userinit\")\n| join kind=inner MasqueradedSpawns on DeviceId\n| where Timestamp between (SpawnTime .. SpawnTime + PersistenceWindow)\n| project SpawnTime, PersistenceTime = Timestamp, DeviceName, AccountName,\n          ParentDoc, MasqFile, MasqPath, MasqSha256, MasqOriginalName, MasqCompany,\n          RegistryKey, RegistryValueName, RegistryValueData,\n          PersistDelaySec = datetime_diff('second', Timestamp, SpawnTime)\n| order by SpawnTime desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet PersistenceWindow = 15m;\nlet OfficeDecoyParents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"hwp.exe\",\"hword.exe\",\"hcell.exe\",\"hshow.exe\",\"acrord32.exe\",\"acrobat.exe\",\"eqnedt32.exe\"]);\nlet MasqueradeNames = dynamic([\"svchost.exe\",\"igfxtray.exe\",\"conhost.exe\",\"wuauclt.exe\",\"ctfmon.exe\",\"taskhost.exe\",\"taskhostw.exe\",\"explorer.exe\",\"lsass.exe\",\"services.exe\",\"winlogon.exe\",\"smss.exe\",\"spoolsv.exe\"]);\nlet MasqueradedSpawns = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ParentName = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1])),\n             ChildName  = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n    | where ParentName in (OfficeDecoyParents)\n    | where ChildName  in (MasqueradeNames)\n    | where not(NewProcessName startswith @\"C:\\Windows\\System32\\\")\n    | where not(NewProcessName startswith @\"C:\\Windows\\SysWOW64\\\")\n    | where not(NewProcessName startswith @\"C:\\Windows\\WinSxS\\\")\n    | where not(NewProcessName startswith @\"C:\\Program Files\\\")\n    | where not(NewProcessName startswith @\"C:\\Program Files (x86)\\\")\n    | project SpawnTime = TimeGenerated, Computer, Account,\n              ParentDoc = ParentName, MasqImage = NewProcessName,\n              MasqCmd = CommandLine, MasqPid = NewProcessId;\nlet RunKeyWrites = ImRegistryEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n    | where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\", @\"\\Explorer\\Run\", @\"\\Userinit\")\n    | project RegTime = TimeGenerated, DvcHostname, RegistryKey, RegistryValue, RegistryValueData, ActingProcessName;\nMasqueradedSpawns\n| join kind=inner RunKeyWrites on $left.Computer == $right.DvcHostname\n| where RegTime between (SpawnTime .. SpawnTime + PersistenceWindow)\n| project SpawnTime, RegTime, Computer, Account,\n          ParentDoc, MasqImage, MasqCmd,\n          RegistryKey, RegistryValue, RegistryValueData, ActingProcessName,\n          PersistDelaySec = datetime_diff('second', RegTime, SpawnTime)\n| order by SpawnTime desc",
      "sigma_yaml": "title: Darkhotel - System-Binary Masquerade Spawned by Office/HWP Decoy\nid: 8a1d3c12-5b9e-4f1a-9c33-darkhotel-masq01\nstatus: experimental\ndescription: Detects an executable named like a legitimate Windows system binary (svchost.exe, igfxtray.exe, etc.) running outside System32/SysWOW64/Program Files, spawned by an Office or Hancom HWP process. Aligned with Darkhotel (G0012) loader masquerading.\nreferences:\n  - https://attack.mitre.org/groups/G0012/\n  - https://securelist.com/the-darkhotel-apt/66779/\nauthor: SOC Detection Engineering\ndate: 2026/05/04\ntags:\n  - attack.initial_access\n  - attack.execution\n  - attack.defense_evasion\n  - attack.t1036.005\n  - attack.t1203\n  - attack.t1566.001\n  - attack.g0012\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent:\n    ParentImage|endswith:\n      - '\\winword.exe'\n      - '\\excel.exe'\n      - '\\powerpnt.exe'\n      - '\\outlook.exe'\n      - '\\hwp.exe'\n      - '\\hword.exe'\n      - '\\hcell.exe'\n      - '\\hshow.exe'\n      - '\\acrord32.exe'\n      - '\\acrobat.exe'\n      - '\\eqnedt32.exe'\n  selection_child_name:\n    Image|endswith:\n      - '\\svchost.exe'\n      - '\\igfxtray.exe'\n      - '\\conhost.exe'\n      - '\\wuauclt.exe'\n      - '\\ctfmon.exe'\n      - '\\taskhost.exe'\n      - '\\taskhostw.exe'\n      - '\\explorer.exe'\n      - '\\lsass.exe'\n      - '\\services.exe'\n      - '\\winlogon.exe'\n      - '\\smss.exe'\n      - '\\spoolsv.exe'\n  filter_legit_paths:\n    Image|startswith:\n      - 'C:\\Windows\\System32\\'\n      - 'C:\\Windows\\SysWOW64\\'\n      - 'C:\\Windows\\WinSxS\\'\n      - 'C:\\Program Files\\'\n      - 'C:\\Program Files (x86)\\'\n  condition: selection_parent and selection_child_name and not filter_legit_paths\nfalsepositives:\n  - Some legacy installers drop helper binaries with these names; tune by image hash signer.\nlevel: high",
      "rationale": "Darkhotel \u2014 flagged as DUBNIUM by Microsoft and Zigzag Hail by Mandiant \u2014 has been observed delivering loaders via Korean Hancom HWP attachments and weaponised Office/PDF docs (T1566.001 + T1203), then writing a renamed payload to user-writable directories that mimics a Windows system binary (T1036.005) and chaining HKCU\\...\\Run for low-privilege persistence (T1547.001). The combined signal \u2014 Offi",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    }
  ]
}