# intel/.llm_actor_uc_cache/d6/d60d3b55edbb5a4a8eb5d4cbbec9b5bc0540db10.json
{
  "name": "Water Galura",
  "ucs": [
    {
      "title": "[LLM] Water Galura / Qilin (Agenda) Safe-Mode-reboot ransomware staging chain",
      "description": "Hunts Qilin's signature pre-encryption sequence: affiliates use `bcdedit /set safeboot` (often `network`) plus shadow-copy deletion and event-log clearing on the same host within minutes, then reboot into Safe Mode where EDR drivers don't load. This staging chain is documented by Trend Micro, Sophos and Group-IB as characteristic Qilin / Agenda RaaS tradecraft (Water Galura = G1050 = GOLD FEATHER).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1486",
        "T1490",
        "T1562.009",
        "T1070.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdlines values(Processes.process_name) as binaries values(Processes.parent_process_name) as parents from datamodel=Endpoint.Processes where (Processes.user!=\"*$\" AND Processes.user!=\"SYSTEM\") AND ((Processes.process_name=\"bcdedit.exe\" AND (Processes.process=\"*safeboot*network*\" OR Processes.process=\"*safeboot*minimal*\" OR Processes.process=\"*bootstatuspolicy*ignoreallfailures*\" OR Processes.process=\"*recoveryenabled*no*\")) OR (Processes.process_name=\"vssadmin.exe\" AND Processes.process=\"*delete*shadows*\") OR (Processes.process_name=\"wmic.exe\" AND Processes.process=\"*shadowcopy*delete*\") OR (Processes.process_name=\"wbadmin.exe\" AND Processes.process=\"*delete*catalog*\") OR (Processes.process_name=\"wevtutil.exe\" AND Processes.process=\"*cl*\") OR (Processes.process_name=\"powershell*\" AND (Processes.process=\"*Remove-WmiObject*Win32_Shadowcopy*\" OR Processes.process=\"*Set-MpPreference*Disable*\"))) by _time span=1h Processes.dest Processes.user | `drop_dm_object_name(Processes)` | eval qilin_chain=if(match(mvjoin(binaries,\",\"),\"bcdedit\\.exe\") AND (match(mvjoin(binaries,\",\"),\"vssadmin\\.exe\") OR match(mvjoin(binaries,\",\"),\"wmic\\.exe\") OR match(mvjoin(binaries,\",\"),\"wbadmin\\.exe\") OR match(mvjoin(binaries,\",\"),\"wevtutil\\.exe\")),1,0) | where qilin_chain=1 | sort - lastTime",
      "kql": "let LookbackDays = 7d;\nlet WindowMinutes = 30m;\n// Step 1 \u2014 Qilin's pathognomonic move: switch the boot configuration to Safe Mode (often `network` so SMB shares stay reachable for the encryptor)\nlet SafeBoot = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where FileName =~ \"bcdedit.exe\"\n    | where ProcessCommandLine has \"safeboot\"\n       or (ProcessCommandLine has \"bootstatuspolicy\" and ProcessCommandLine has \"ignoreallfailures\")\n    | project SafeBootTime = Timestamp, DeviceId, DeviceName, AccountName,\n              SafeBootCmd = ProcessCommandLine,\n              SafeBootParent = InitiatingProcessFileName;\n// Step 2 \u2014 Inhibit recovery (VSS) and/or clear logs on the same host within +/- 30 minutes\nlet Wreckage = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where (FileName =~ \"vssadmin.exe\" and ProcessCommandLine has \"delete\" and ProcessCommandLine has \"shadows\")\n        or (FileName =~ \"wmic.exe\" and ProcessCommandLine has_all (\"shadowcopy\",\"delete\"))\n        or (FileName =~ \"wbadmin.exe\" and ProcessCommandLine has \"delete\" and ProcessCommandLine has_any (\"catalog\",\"backup\"))\n        or (FileName =~ \"wevtutil.exe\" and ProcessCommandLine has \"cl\")\n        or (FileName in~ (\"powershell.exe\",\"pwsh.exe\")\n            and ProcessCommandLine has_any (\"Remove-WmiObject Win32_Shadowcopy\",\n                                            \"Get-WmiObject Win32_Shadowcopy\",\n                                            \"Set-MpPreference -DisableRealtimeMonitoring\"))\n    | project WreckTime = Timestamp, DeviceId, WreckCmd = ProcessCommandLine, WreckBin = FileName;\nSafeBoot\n| join kind=inner Wreckage on DeviceId\n| where WreckTime between ((SafeBootTime - WindowMinutes) .. (SafeBootTime + WindowMinutes))\n| summarize SafeBootCmd = any(SafeBootCmd),\n            DestructiveCmds = make_set(WreckCmd, 10),\n            DestructiveBins = make_set(WreckBin, 10),\n            FirstStep = min(SafeBootTime),\n            LastStep  = max(WreckTime)\n            by DeviceName, AccountName\n| order by FirstStep desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet WindowMinutes = 30m;\nlet SafeBoot = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where AccountType == \"User\"\n    | where Account !endswith \"$\"\n    | where NewProcessName endswith \"\\\\bcdedit.exe\"\n    | where CommandLine has \"safeboot\"\n       or (CommandLine has \"bootstatuspolicy\" and CommandLine has \"ignoreallfailures\")\n    | project SafeBootTime = TimeGenerated, Computer, Account,\n              SafeBootCmd = CommandLine,\n              SafeBootParent = ParentProcessName;\nlet Wreckage = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | where (NewProcessName endswith \"\\\\vssadmin.exe\" and CommandLine has \"delete\" and CommandLine has \"shadows\")\n        or (NewProcessName endswith \"\\\\wmic.exe\" and CommandLine has \"shadowcopy\" and CommandLine has \"delete\")\n        or (NewProcessName endswith \"\\\\wbadmin.exe\" and CommandLine has \"delete\" and (CommandLine has \"catalog\" or CommandLine has \"backup\"))\n        or (NewProcessName endswith \"\\\\wevtutil.exe\" and CommandLine has \"cl\")\n        or (NewProcessName has_any (\"powershell.exe\",\"pwsh.exe\")\n            and (CommandLine has \"Remove-WmiObject Win32_Shadowcopy\"\n              or CommandLine has \"Get-WmiObject Win32_Shadowcopy\"\n              or CommandLine has \"Set-MpPreference -DisableRealtimeMonitoring\"))\n    | project WreckTime = TimeGenerated, Computer, WreckCmd = CommandLine, WreckBin = NewProcessName;\nSafeBoot\n| join kind=inner Wreckage on Computer\n| where WreckTime between ((SafeBootTime - WindowMinutes) .. (SafeBootTime + WindowMinutes))\n| summarize SafeBootCmd = any(SafeBootCmd),\n            DestructiveCmds = make_set(WreckCmd, 10),\n            DestructiveBins = make_set(WreckBin, 10),\n            FirstStep = min(SafeBootTime),\n            LastStep  = max(WreckTime)\n            by Computer, Account\n| order by FirstStep desc",
      "sigma_yaml": "",
      "rationale": "Qilin/Agenda affiliates (Water Galura) are documented by Sophos, Trend Micro and Group-IB as one of the few RaaS crews that consistently boot endpoints into Safe Mode (often `safeboot network`) before running the encryptor \u2014 a deliberate EDR-evasion step that Conti, BlackCat and others rarely use. Pairing the bcdedit pivot with VSS deletion/log clearing within a 30-minute window on the same host i"
    },
    {
      "title": "[LLM] Qilin (Water Galura) GPO-deployed Chrome credential harvest (IPScanner.ps1 pattern)",
      "description": "Hunts the novel Qilin-affiliate tradecraft documented by Sophos X-Ops in August 2024: after taking the domain controller, the affiliate pushed a Group Policy that ran a logon scheduled task launching a PowerShell script (referenced internally as `IPScanner.ps1` / `logon.ps1`) that copied Chrome's `Login Data` and `Local State` files from every domain endpoint to a SYSVOL or share-resident output. This is a uniquely Qilin TTP \u2014 credential harvesting at domain scale via GPO is rarely observed.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1555.003",
        "T1484.001",
        "T1059.001",
        "T1053.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdlines values(Processes.parent_process_name) as parents values(Processes.user) as users from datamodel=Endpoint.Processes where Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\") AND Processes.parent_process_name IN (\"svchost.exe\",\"taskeng.exe\",\"schtasks.exe\",\"explorer.exe\",\"userinit.exe\") AND (Processes.process=\"*\\\\Google\\\\Chrome\\\\User Data\\\\*\" OR Processes.process=\"*Login Data*\" OR Processes.process=\"*Local State*\" OR Processes.process=\"*\\\\Microsoft\\\\Edge\\\\User Data\\\\*\" OR Processes.process=\"*IPScanner*\" OR Processes.process=\"*\\\\SYSVOL\\\\*.ps1*\") by Processes.dest Processes.user Processes.process Processes.parent_process_name _time span=10m | `drop_dm_object_name(Processes)` | stats dc(dest) as host_count values(dest) as hosts values(process) as cmdlines values(parent_process_name) as parents min(firstTime) as firstSeen max(lastTime) as lastSeen by user | where host_count >= 3 | sort - host_count",
      "kql": "let LookbackDays = 7d;\nlet BrowserCredArtifacts = dynamic([@\"\\Google\\Chrome\\User Data\\\", @\"\\Microsoft\\Edge\\User Data\\\", \"Login Data\", \"Local State\", \"Web Data\", \"Cookies\"]);\n// Hit 1 \u2014 PowerShell launched from a GPO/scheduled-task lineage AND touching Chrome cred artefacts\nlet ProcessSignal = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where FileName in~ (\"powershell.exe\",\"pwsh.exe\")\n    | where InitiatingProcessFileName in~ (\"svchost.exe\",\"taskeng.exe\",\"schtasks.exe\",\"userinit.exe\",\"explorer.exe\",\"gpscript.exe\")\n    | where ProcessCommandLine has_any (@\"Chrome\\User Data\", @\"Edge\\User Data\", \"Login Data\", \"Local State\", \"IPScanner\", \"logon.ps1\", @\"\\SYSVOL\\\")\n       or InitiatingProcessCommandLine has @\"\\SYSVOL\\\"\n    | project Timestamp, DeviceId, DeviceName, AccountName,\n              ProcessCommandLine, InitiatingProcessFileName, InitiatingProcessCommandLine;\n// Hit 2 \u2014 file-system corroboration: PowerShell creating/copying the actual cred files\nlet FileSignal = DeviceFileEvents\n    | where Timestamp > ago(LookbackDays)\n    | where ActionType in (\"FileCreated\",\"FileRenamed\",\"FileModified\")\n    | where InitiatingProcessFileName in~ (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\")\n    | where FileName has_any (\"Login Data\",\"LoginData\",\"Local State\",\"LocalState\",\"Cookies\",\"Web Data\")\n        or PreviousFileName has_any (\"Login Data\",\"Local State\",\"Cookies\")\n    | project Timestamp, DeviceId, DeviceName, FileName, FolderPath,\n              InitiatingProcessCommandLine, InitiatingProcessAccountName;\n// Anti-baseline \u2014 fan-out across multiple devices in a short window is the GPO-push tell\nProcessSignal\n| union FileSignal\n| summarize Hosts = dcount(DeviceId),\n            HostList = make_set(DeviceName, 25),\n            Cmds = make_set(ProcessCommandLine, 10),\n            FirstSeen = min(Timestamp),\n            LastSeen  = max(Timestamp)\n            by bin(Timestamp, 30m), AccountName\n| where Hosts >= 3   // GPO push hits many endpoints near-simultaneously; tune for org size\n| order by FirstSeen desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet ProcSignal = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where NewProcessName has_any (\"\\\\powershell.exe\",\"\\\\pwsh.exe\")\n    | where ParentProcessName has_any (\"\\\\svchost.exe\",\"\\\\taskeng.exe\",\"\\\\schtasks.exe\",\"\\\\userinit.exe\",\"\\\\explorer.exe\",\"\\\\gpscript.exe\")\n    | where CommandLine has_any (@\"Chrome\\User Data\", @\"Edge\\User Data\", \"Login Data\", \"Local State\", \"IPScanner\", \"logon.ps1\", \"\\\\SYSVOL\\\\\")\n    | project TimeGenerated, Computer, Account, CommandLine, ParentProcessName;\nlet GpoSignal = AuditLogs\n    | where TimeGenerated > ago(LookbackDays)\n    | where OperationName has_any (\"Update group policy\", \"Set group policy\")\n    | extend Initiator = tostring(parse_json(tostring(InitiatedBy)).user.userPrincipalName),\n             TargetGpo = tostring(TargetResources[0].displayName)\n    | project TimeGenerated, Initiator, TargetGpo, OperationName;\nProcSignal\n| summarize Hosts = dcount(Computer),\n            HostList = make_set(Computer, 25),\n            Cmds = make_set(CommandLine, 10),\n            FirstSeen = min(TimeGenerated),\n            LastSeen  = max(TimeGenerated)\n            by bin(TimeGenerated, 30m), Account\n| where Hosts >= 3   // simultaneous fan-out across endpoints == GPO-driven push\n| order by FirstSeen desc",
      "sigma_yaml": "title: Qilin/Water Galura \u2014 PowerShell touching Chrome credential store under scheduled-task lineage\nid: 7c3f1e6d-9d6a-4f2b-9b3e-9b9c8a1f7e21\nstatus: experimental\ndescription: Detects the Qilin/Agenda affiliate TTP of harvesting Chrome Login Data and Local State across the domain via a GPO-pushed logon scheduled task launching PowerShell.\nreferences:\n  - https://news.sophos.com/en-us/2024/08/22/qilin-ransomware-caught-stealing-credentials-stored-in-google-chrome\n  - https://attack.mitre.org/groups/G1050/\nauthor: Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.credential_access\n  - attack.t1555.003\n  - attack.t1484.001\n  - attack.execution\n  - attack.t1059.001\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_proc:\n    Image|endswith:\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n    ParentImage|endswith:\n      - '\\svchost.exe'\n      - '\\taskeng.exe'\n      - '\\schtasks.exe'\n      - '\\userinit.exe'\n      - '\\gpscript.exe'\n  selection_cmd:\n    CommandLine|contains:\n      - 'Chrome\\User Data'\n      - 'Edge\\User Data'\n      - 'Login Data'\n      - 'Local State'\n      - 'IPScanner'\n      - 'logon.ps1'\n      - '\\SYSVOL\\'\n  condition: selection_proc and selection_cmd\nfalsepositives:\n  - Legitimate enterprise password-migration scripts deployed via GPO (rare; should be allow-listed by hash/path)\nlevel: high\n",
      "rationale": "Sophos X-Ops attributed this exact TTP \u2014 GPO-pushed PowerShell logon task copying Chrome's `Login Data` and `Local State` from every domain endpoint to a share \u2014 to a Qilin (Water Galura) affiliate in mid-2024 and labelled it 'a worrying evolution' precisely because it's not seen elsewhere. The GPO-push delivery model produces a near-simultaneous fan-out signature across many hosts at once, which "
    }
  ]
}