# intel/.llm_actor_uc_cache/88/883c9c31729fd6381b78ba85a02f2b213a7a2a3a.json
{
  "name": "BRONZE BUTLER",
  "ucs": [
    {
      "title": "[LLM] BRONZE BUTLER (Tick) \u2014 Daserf/xxmm DLL side-loading via signed-binary proxy in user-writable paths",
      "description": "Hunts BRONZE BUTLER's signature loader chain: a legitimately signed third-party binary (commonly a Symantec, Trend Micro, Hewlett-Packard, or Kaspersky utility) is dropped into AppData/ProgramData/Public/Temp and used to sideload a co-resident malicious DLL (Daserf, xxmm/Minzen, Datper). This T1574.001 + T1036.005 + T1027 chain is one of Tick's most consistent fingerprints across a decade of intrusions in Japan.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.001",
        "T1036.005",
        "T1027.001",
        "T1071.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as command_line values(Processes.parent_process_name) as parent values(Processes.process_hash) as proc_hash from datamodel=Endpoint.Processes where Processes.process_path IN (\"*\\\\AppData\\\\Local\\\\*\",\"*\\\\AppData\\\\Roaming\\\\*\",\"*\\\\ProgramData\\\\*\",\"*\\\\Users\\\\Public\\\\*\",\"*\\\\Windows\\\\Temp\\\\*\",\"*\\\\Downloads\\\\*\") Processes.process_name IN (\"rcsync.exe\",\"vncviewer.exe\",\"GUP.exe\",\"avp.exe\",\"avpui.exe\",\"DefenderUI.exe\",\"PresentationHost.exe\",\"hpqtra08.exe\",\"hkcmd.exe\",\"smax4pnp.exe\",\"NMSAccessU.exe\",\"OneDrive.exe\",\"ccSvcHst.exe\",\"SMC.exe\",\"NavLogon.exe\") by Processes.dest Processes.user Processes.process_name Processes.process_path Processes.process_id _time\n| `drop_dm_object_name(Processes)`\n| eval proc_dir=replace(process_path,\"\\\\\\\\[^\\\\\\\\]+$\",\"\")\n| join type=inner dest proc_dir [\n    | tstats `security_content_summariesonly` count values(Filesystem.file_name) as dropped_dll values(Filesystem.process_name) as dropper values(Filesystem.file_hash) as dll_hash from datamodel=Endpoint.Filesystem where Filesystem.action=created Filesystem.file_name=\"*.dll\" by Filesystem.dest Filesystem.file_path\n    | `drop_dm_object_name(Filesystem)`\n    | eval proc_dir=replace(file_path,\"\\\\\\\\[^\\\\\\\\]+$\",\"\")\n    | fields dest proc_dir dropped_dll dropper dll_hash\n  ]\n| where dropper!=process_name\n| `security_content_ctime(firstTime)`\n| `security_content_ctime(lastTime)`\n| table firstTime lastTime dest user process_name process_path proc_hash command_line dropped_dll dll_hash dropper parent",
      "kql": "// BRONZE BUTLER / Tick \u2014 signed legit binary side-loading from user-writable path\nlet _suspicious_paths = dynamic([@\"\\\\AppData\\\\\", @\"\\\\ProgramData\\\\\", @\"\\\\Users\\\\Public\\\\\", @\"\\\\Windows\\\\Temp\\\\\", @\"\\\\Downloads\\\\\"]);\nlet _vendor_decoys = dynamic([\"Symantec\",\"Trend Micro\",\"Kaspersky\",\"Hewlett-Packard\",\"HP Inc.\",\"Hewlett Packard\",\"Sourcenext\",\"Lhaplus\",\"NEC\",\"Fujitsu\",\"Justsystems\",\"Sky Co.,LTD.\"]);\nlet SignedBinaryRunFromUserPath =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where AccountName !endswith \"$\"\n    | where InitiatingProcessVersionInfoCompanyName has_any (_vendor_decoys)\n         or ProcessVersionInfoCompanyName has_any (_vendor_decoys)\n    | where FolderPath has_any (_suspicious_paths)\n    | where not(FolderPath startswith @\"C:\\Program Files\") and not(FolderPath startswith @\"C:\\Program Files (x86)\")\n    | project ProcTime = Timestamp, DeviceId, DeviceName, AccountName,\n              FileName, FolderPath, SHA256, ProcessCommandLine,\n              ParentImage = InitiatingProcessFileName,\n              ParentCmd   = InitiatingProcessCommandLine,\n              ProcCompany = ProcessVersionInfoCompanyName,\n              ParentCompany = InitiatingProcessVersionInfoCompanyName;\n// Side-loaded DLL = different SHA256, same folder, loaded by the signed proxy\nDeviceImageLoadEvents\n| where Timestamp > ago(7d)\n| where FileName endswith \".dll\"\n| where InitiatingProcessFolderPath has_any (_suspicious_paths)\n| where FolderPath =~ InitiatingProcessFolderPath           // co-resident DLL drop\n| where SHA256 != InitiatingProcessSHA256\n| where InitiatingProcessVersionInfoCompanyName has_any (_vendor_decoys)\n| join kind=inner SignedBinaryRunFromUserPath on $left.DeviceId == $right.DeviceId, $left.InitiatingProcessSHA256 == $right.SHA256\n| where Timestamp between (ProcTime .. ProcTime + 5m)\n| project LoadTime = Timestamp, DeviceName, AccountName,\n          ProxyBinary = InitiatingProcessFileName,\n          ProxyVendor = InitiatingProcessVersionInfoCompanyName,\n          ProxyPath = InitiatingProcessFolderPath,\n          ProxySHA256 = InitiatingProcessSHA256,\n          SideLoadedDLL = FileName,\n          DLLPath = FolderPath,\n          DLLSHA256 = SHA256,\n          ParentImage, ParentCmd\n| order by LoadTime desc",
      "sentinel_kql": "// BRONZE BUTLER / Tick \u2014 signed legit binary running from a user-writable path is the\n// process-create half of the Daserf side-load chain. Sentinel SecurityEvent does not carry\n// DLL ImageLoad telemetry natively, so we approximate via 4688 + ImFileEvent join.\nlet _suspicious_paths = dynamic([@\"\\AppData\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\Windows\\Temp\\\", @\"\\Downloads\\\"]);\nlet _vendor_proxy_names = dynamic([\"rcsync.exe\",\"vncviewer.exe\",\"GUP.exe\",\"avp.exe\",\"avpui.exe\",\"hpqtra08.exe\",\"NMSAccessU.exe\",\"ccSvcHst.exe\",\"SMC.exe\",\"NavLogon.exe\",\"hkcmd.exe\",\"PresentationHost.exe\"]);\nlet ProxyExec =\n    SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend BinName = tostring(split(NewProcessName, \"\\\\\")[-1])\n    | where BinName in~ (_vendor_proxy_names)\n    | where NewProcessName has_any (_suspicious_paths)\n    | extend ProxyDir = strcat_array(array_slice(split(NewProcessName,\"\\\\\"), 0, -2), \"\\\\\")\n    | project ProxyTime = TimeGenerated, Computer, Account, BinName, NewProcessName, ProxyDir, ParentProcessName;\nImFileEvent\n| where TimeGenerated > ago(7d)\n| where EventType == \"FileCreated\"\n| where TargetFileName endswith \".dll\"\n| where TargetFilePath has_any (_suspicious_paths)\n| extend DropDir = strcat_array(array_slice(split(TargetFilePath,\"\\\\\"), 0, -2), \"\\\\\")\n| join kind=inner ProxyExec on $left.DvcHostname == $right.Computer, $left.DropDir == $right.ProxyDir\n| where TimeGenerated between (ProxyTime - 1h .. ProxyTime + 5m)\n| project DropTime = TimeGenerated, ProxyTime, Computer, Account,\n          ProxyBinary = BinName, ProxyPath = NewProcessName,\n          SideLoadedDLL = TargetFileName, DLLPath = TargetFilePath, DLLSHA256 = TargetFileSHA256,\n          ParentProcessName\n| order by DropTime desc",
      "sigma_yaml": "",
      "rationale": "Tick / BRONZE BUTLER (G0060) has used the same loader recipe across Daserf, xxmm/Minzen, Datper and Ghostdown campaigns \u2014 drop a legitimately signed Japanese-market or AV vendor binary alongside a malicious DLL into a user-writable directory and let DLL search-order hijacking run the implant under a trusted parent. Trend Micro (Nov 2017) and Secureworks (Oct 2017) both fingerprint this chain. The "
    },
    {
      "title": "[LLM] BRONZE BUTLER (Tick) \u2014 Office/archive-spawned 'at.exe' or schtasks persistence (Daserf installer pattern)",
      "description": "Catches Tick's documented post-exploitation step: after a spearphishing LZH/RAR/maldoc payload executes, the unpacked dropper registers persistence with the legacy 'at.exe' command or schtasks /create \u2014 frequently invoking a binary in AppData/Public/Temp. The Office/archive-handler parent + at.exe child combination is highly characteristic of Tick installers (T1566.001 \u2192 T1204.002 \u2192 T1053.002/T1053.005).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "medium",
      "techs": [
        "T1053.002",
        "T1053.005",
        "T1204.002",
        "T1566.001",
        "T1059.003"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.process_path) as path values(Processes.parent_process) as parent_cmdline from datamodel=Endpoint.Processes where Processes.process_name IN (\"at.exe\",\"schtasks.exe\") (Processes.process=\"*create*\" OR Processes.process=\"* /sc *\" OR Processes.process=\"* /ru *\" OR Processes.process=\"* /tr *\") Processes.parent_process_name IN (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"explorer.exe\",\"WinRAR.exe\",\"Lhaplus.exe\",\"7zFM.exe\",\"7zG.exe\") by Processes.dest Processes.user Processes.process_name Processes.parent_process_name\n| `drop_dm_object_name(Processes)`\n| where match(cmdline,\"(?i)\\\\\\\\AppData\\\\\\\\|\\\\\\\\ProgramData\\\\\\\\|\\\\\\\\Users\\\\\\\\Public\\\\\\\\|\\\\\\\\Temp\\\\\\\\|\\\\\\\\Downloads\\\\\\\\\")\n| `security_content_ctime(firstTime)`\n| `security_content_ctime(lastTime)`\n| table firstTime lastTime dest user parent_process_name parent_cmdline process_name cmdline path",
      "kql": "// BRONZE BUTLER \u2014 at.exe / schtasks persistence chained from Office or archive-handler parent\nlet _office_archive_parents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"explorer.exe\",\"winrar.exe\",\"lhaplus.exe\",\"7zfm.exe\",\"7zg.exe\",\"hh.exe\"]);\nlet _user_writable = dynamic([@\"\\\\AppData\\\\\", @\"\\\\ProgramData\\\\\", @\"\\\\Users\\\\Public\\\\\", @\"\\\\Windows\\\\Temp\\\\\", @\"\\\\Downloads\\\\\"]);\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where AccountName !endswith \"$\"\n| where FileName in~ (\"at.exe\",\"schtasks.exe\")\n| where InitiatingProcessFileName in~ (_office_archive_parents)\n   or InitiatingProcessParentFileName in~ (_office_archive_parents)\n| where ProcessCommandLine has_any (\"/create\",\"/sc \",\"/tr \",\"/ru \") or FileName =~ \"at.exe\"\n| where ProcessCommandLine has_any (_user_writable)\n   or ProcessCommandLine matches regex @\"(?i)\\\\(AppData|ProgramData|Public|Temp|Downloads)\\\\\"\n| project Timestamp, DeviceName, AccountName,\n          Tool = FileName,\n          TaskCmd = ProcessCommandLine,\n          ParentImage = InitiatingProcessFileName,\n          ParentCmd   = InitiatingProcessCommandLine,\n          GrandParent = InitiatingProcessParentFileName,\n          IntegrityLevel = ProcessIntegrityLevel,\n          IsRemote = AdditionalFields\n| order by Timestamp desc",
      "sentinel_kql": "// BRONZE BUTLER \u2014 at.exe / schtasks persistence chained from Office or archive parent\nlet _office_archive_parents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"explorer.exe\",\"winrar.exe\",\"lhaplus.exe\",\"7zfm.exe\",\"7zg.exe\",\"hh.exe\"]);\nlet _user_writable = dynamic([@\"\\AppData\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\Windows\\Temp\\\", @\"\\Downloads\\\"]);\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| extend ChildBin = tolower(tostring(split(NewProcessName, \"\\\\\")[-1])),\n         ParentBin = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1]))\n| where ChildBin in (\"at.exe\",\"schtasks.exe\")\n| where ParentBin in (_office_archive_parents)\n| where CommandLine has_any (\"/create\",\"/sc \",\"/tr \",\"/ru \") or ChildBin == \"at.exe\"\n| where CommandLine has_any (_user_writable)\n   or CommandLine matches regex @\"(?i)\\\\(AppData|ProgramData|Public|Temp|Downloads)\\\\\"\n| project TimeGenerated, Computer, Account, Tool = ChildBin, TaskCmd = CommandLine,\n          ParentImage = ParentBin, ParentProcessName, TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "title: BRONZE BUTLER Persistence via at.exe or schtasks Spawned by Office or Archive Handler\nid: 5b9c2d3e-7c1b-4d8e-9c10-tick-at-schtasks\nstatus: experimental\ndescription: Detects Tick / BRONZE BUTLER post-exploitation persistence where at.exe or schtasks.exe is invoked by an Office product or Japanese archive handler (Lhaplus / WinRAR / 7-Zip) and registers a task pointing at a user-writable path.\nreferences:\n  - https://attack.mitre.org/groups/G0060/\n  - https://www.secureworks.com/research/bronze-butler-targets-japanese-businesses\n  - https://documents.trendmicro.com/assets/Tech-Brief-Untangling-the-Patchwork-Cyberespionage-Group.pdf\nauthor: SOC Detection Engineering\ndate: 2026/05/05\ntags:\n  - attack.persistence\n  - attack.t1053.002\n  - attack.t1053.005\n  - attack.t1204.002\n  - attack.g0060\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_child:\n    Image|endswith:\n      - '\\at.exe'\n      - '\\schtasks.exe'\n  selection_parent:\n    ParentImage|endswith:\n      - '\\winword.exe'\n      - '\\excel.exe'\n      - '\\powerpnt.exe'\n      - '\\outlook.exe'\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n      - '\\mshta.exe'\n      - '\\rundll32.exe'\n      - '\\explorer.exe'\n      - '\\winrar.exe'\n      - '\\Lhaplus.exe'\n      - '\\7zFM.exe'\n      - '\\7zG.exe'\n      - '\\hh.exe'\n  selection_writable:\n    CommandLine|contains:\n      - '\\AppData\\'\n      - '\\ProgramData\\'\n      - '\\Users\\Public\\'\n      - '\\Windows\\Temp\\'\n      - '\\Downloads\\'\n  condition: selection_child and selection_parent and selection_writable\nfields:\n  - User\n  - ParentImage\n  - ParentCommandLine\n  - Image\n  - CommandLine\nfalsepositives:\n  - Legitimate IT scripts deploying tasks via Office macros (rare; allowlist by signed publisher)\nlevel: high",
      "rationale": "Tick has historically reached for the legacy 'at' command \u2014 unusual on modern Windows and almost never invoked by Office or archive handlers in legitimate operations \u2014 alongside schtasks /create to anchor Daserf and follow-on tools. The parent set explicitly includes Lhaplus, the Japanese-market archiver Tick weaponises in spearphishing LZH bundles per Trend Micro 2017/2019 reporting. This combina",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    }
  ]
}