# intel/.llm_actor_uc_cache/ea/ea8bae7276f1bf64b4a98aacb0877169ae68121d.json
{
  "name": "Carbanak",
  "ucs": [
    {
      "title": "[LLM] Carbanak/FIN7 service persistence: rundll32-loaded DLL registered as masqueraded Windows service",
      "description": "Hunts Carbanak's signature install chain where the implant DLL is dropped (commonly under %ProgramData%, %AppData%, or C:\\Windows\\Help\\) and registered as a Windows service whose DisplayName masquerades as a legitimate Microsoft component (e.g. 'Microsoft Help', 'Windows Defender Helper', 'WmiPrvSE') with svchost/rundll32 hosting the payload. This combines T1036.004 + T1218.011 + T1543.003 in the exact pattern Kaspersky and FireEye documented for the Anunak/Carbanak backdoor.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1543.003",
        "T1036.004",
        "T1036.005",
        "T1218.011"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.process_name=\"rundll32.exe\" (Processes.process=\"*ProgramData*\\\\*.dll*\" OR Processes.process=\"*AppData*\\\\*.dll*\" OR Processes.process=\"*\\\\Windows\\\\Help\\\\*.dll*\" OR Processes.process=\"*\\\\Windows\\\\Tasks\\\\*.dll*\") by Processes.dest Processes.user Processes.parent_process_name Processes.process Processes.process_guid | `drop_dm_object_name(Processes)` | join type=inner process_guid [| tstats summariesonly=true count from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\CurrentControlSet\\\\Services\\\\*\" (Registry.registry_value_name=\"ImagePath\" OR Registry.registry_value_name=\"ServiceDll\" OR Registry.registry_value_name=\"DisplayName\") (Registry.registry_value_data=\"*Microsoft Help*\" OR Registry.registry_value_data=\"*Windows Defender Helper*\" OR Registry.registry_value_data=\"*WmiPrvSE*\" OR Registry.registry_value_data=\"*rundll32*ProgramData*\" OR Registry.registry_value_data=\"*rundll32*AppData*\") by Registry.dest Registry.process_guid Registry.registry_value_name Registry.registry_value_data | `drop_dm_object_name(Registry)`] | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` | table firstTime lastTime dest user parent_process_name process registry_value_name registry_value_data",
      "kql": "let suspect_paths = dynamic([@\"\\programdata\\\", @\"\\appdata\\\", @\"\\windows\\help\\\", @\"\\windows\\tasks\\\", @\"\\windows\\debug\\\"]);\nlet masq_names = dynamic([\"Microsoft Help\",\"Windows Defender Helper\",\"WmiPrvSE\",\"Windows Update Helper\",\"Microsoft Network Realtime Inspection\"]);\nlet rundll = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName =~ \"rundll32.exe\"\n| where ProcessCommandLine has_any (suspect_paths)\n| where ProcessCommandLine matches regex @\"(?i)\\.dll[, \\x22\\x27]\\s*[A-Za-z0-9_]+\"\n| project rundll_time=Timestamp, DeviceId, DeviceName, AccountName, InitiatingProcessFileName, ProcessCommandLine, RundllPid=ProcessId;\nlet svc = DeviceRegistryEvents\n| where Timestamp > ago(7d)\n| where RegistryKey has @\"\\CurrentControlSet\\Services\\\"\n| where RegistryValueName in~ (\"ImagePath\",\"ServiceDll\",\"DisplayName\")\n| where (RegistryValueData has_any (masq_names)) or (RegistryValueData has \"rundll32\" and RegistryValueData has_any (suspect_paths))\n| project svc_time=Timestamp, DeviceId, RegistryKey, RegistryValueName, RegistryValueData, SvcInitiator=InitiatingProcessFileName;\nrundll\n| join kind=inner svc on DeviceId\n| where abs(datetime_diff('minute', svc_time, rundll_time)) <= 30\n| project rundll_time, svc_time, DeviceName, AccountName, ProcessCommandLine, RegistryKey, RegistryValueName, RegistryValueData, InitiatingProcessFileName, SvcInitiator",
      "rationale": "Generic 'new service' rules drown SOCs in false positives. This UC is keyed to Carbanak's documented combo: implant dropped to writable user-or-help paths, hosted by rundll32 (T1218.011), then re-launched as a service whose name spoofs Microsoft branding (T1036.004/.005) \u2014 exactly how Kaspersky's 2015 Carbanak APT report and FireEye's FIN7 reports describe the persistence step. Requiring the rundl"
    },
    {
      "title": "[LLM] Carbanak hands-on-keyboard: remote-access tool (Ammyy/RAdmin/TeamViewer) launched on finance/SWIFT/ATM-management host",
      "description": "Carbanak/Anunak operators are known for installing legitimate remote-access software (Ammyy Admin, RAdmin, TeamViewer, VNC) on operator workstations inside banks to perform manual SWIFT message and ATM card-management fraud (T1219). This hunt fires when one of those binaries executes on a host tagged as finance/treasury/ATM-mgmt and beacons to a non-corporate destination \u2014 a high-signal pattern unique to financially motivated intrusion sets like Carbanak and its FIN7/Cobalt offshoots.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "high",
      "techs": [
        "T1219",
        "T1078",
        "T1102.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where (Processes.process_name IN (\"AA_v3.exe\",\"AMMYY_Admin.exe\",\"ammyy.exe\",\"r_server.exe\",\"radmin.exe\",\"rserver3.exe\",\"TeamViewer.exe\",\"winvnc.exe\",\"tvnserver.exe\",\"AnyDesk.exe\")) by Processes.dest Processes.user Processes.process_name Processes.process Processes.parent_process_name Processes.process_guid | `drop_dm_object_name(Processes)` | lookup asset_lookup dest OUTPUT category business_unit | search (category=\"finance\" OR category=\"swift\" OR category=\"treasury\" OR category=\"atm-mgmt\" OR business_unit=\"payments\") | join type=inner process_guid [| tstats summariesonly=true count values(All_Traffic.dest_ip) as dest_ips values(All_Traffic.dest_port) as dest_ports from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest_port IN (443,80,5650,5800,5900,4899,3389) NOT (All_Traffic.dest_category=\"corporate\" OR All_Traffic.dest_ip IN (\"10.0.0.0/8\",\"172.16.0.0/12\",\"192.168.0.0/16\")) by All_Traffic.process_guid] | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` | table firstTime lastTime dest category business_unit user process_name process parent_process_name dest_ips dest_ports",
      "kql": "let rat_bins = dynamic([\"AA_v3.exe\",\"AMMYY_Admin.exe\",\"ammyy.exe\",\"r_server.exe\",\"radmin.exe\",\"rserver3.exe\",\"TeamViewer.exe\",\"winvnc.exe\",\"tvnserver.exe\",\"AnyDesk.exe\",\"vncserver.exe\"]);\nlet finance_tags = dynamic([\"finance\",\"swift\",\"treasury\",\"payments\",\"atm-mgmt\",\"card-ops\"]);\nlet finance_hosts = DeviceInfo\n| where Timestamp > ago(1d)\n| where DeviceCategory has_any (finance_tags) or MachineGroup has_any (finance_tags) or Tags has_any (finance_tags)\n| distinct DeviceId, DeviceName;\nlet rat_exec = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName in~ (rat_bins)\n| project rat_time=Timestamp, DeviceId, DeviceName, AccountName, FileName, ProcessCommandLine, ProcessId, InitiatingProcessFileName;\nlet rat_net = DeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (rat_bins)\n| where RemoteIPType !in (\"Private\",\"Loopback\",\"Reserved\")\n| where RemotePort in (443, 80, 5650, 5800, 5900, 4899, 3389, 7070)\n| project net_time=Timestamp, DeviceId, RemoteIP, RemotePort, RemoteUrl, InitiatingProcessFileName, InitiatingProcessId=InitiatingProcessId;\nfinance_hosts\n| join kind=inner rat_exec on DeviceId\n| join kind=inner rat_net on DeviceId\n| where abs(datetime_diff('minute', net_time, rat_time)) <= 60\n| project rat_time, net_time, DeviceName, AccountName, FileName, ProcessCommandLine, RemoteIP, RemotePort, RemoteUrl",
      "rationale": "Carbanak's whole monetisation model \u2014 documented by Kaspersky (2015), Group-IB Anunak report, and Europol's takedown bulletin \u2014 is hands-on-keyboard fraud against SWIFT terminals, card-management consoles, and ATM controllers using off-the-shelf RATs (Ammyy in particular is the Carbanak signature, RAdmin for FIN7). Asset-tagging the alert to finance/payments/ATM-mgmt hosts and requiring an externa"
    }
  ]
}