# use_cases/install/UC_SCHEDULED_TASK.yml
id: UC_SCHEDULED_TASK
title: Scheduled task created with suspicious image / encoded args
kill_chain: install
confidence: High
description: schtasks.exe /create or Microsoft-Windows-TaskScheduler EventID 4698 with LOLBin actions.
implementations:
- splunk
- defender
mitre_attack:
- id: T1053.005
  name: Scheduled Task
data_models:
  splunk:
  - Endpoint.Processes
  defender:
  - DeviceProcessEvents
splunk_spl: "| tstats `summariesonly` count min(_time) as firstTime max(_time) as lastTime\n    from datamodel=Endpoint.Processes\n    where Processes.process_name=\"schtasks.exe\" AND Processes.process=\"*/create*\"\n      AND (Processes.process=\"*powershell*\" OR Processes.process=\"*cmd.exe*\"\n        OR Processes.process=\"*rundll32*\" OR Processes.process=\"*-enc*\"\n        OR Processes.process=\"*FromBase64*\" OR Processes.process=\"*\\Users\\Public*\"\n        OR Processes.process=\"*\\AppData\\*\")\n    by Processes.dest, Processes.user, Processes.process, Processes.parent_process_name\n| `drop_dm_object_name(Processes)`"
defender_kql: 'DeviceProcessEvents

  | where Timestamp > ago(7d)

  | where FileName =~ "schtasks.exe"

  | where ProcessCommandLine has "/create"

  | where ProcessCommandLine has_any ("powershell","cmd.exe","rundll32","-enc","FromBase64","\Users\Public","\AppData\")

  | project Timestamp, DeviceName, AccountName, ProcessCommandLine, InitiatingProcessFileName'
tier: hunting
fp_rate_estimate: high
required_telemetry:
- Windows Event ID 4698
- Sysmon EID 1
- Defender DeviceProcessEvents
