# intel/.llm_actor_uc_cache/45/45d3450ab2884a5cf508e69da9216e74e5531b4c.json
{
  "name": "APT3",
  "ucs": [
    {
      "title": "[LLM] APT3 (Gothic Panda/Buckeye) accessibility-features RDP backdoor via IFEO Debugger hijack",
      "description": "Hunts the signature APT3/Buckeye persistence trick of redirecting Windows accessibility binaries (sethc.exe sticky-keys, utilman.exe, osk.exe, narrator.exe, magnify.exe, displayswitch.exe) to cmd/powershell via the Image File Execution Options 'Debugger' value or by replacing the binary on disk \u2014 granting SYSTEM shell at the pre-auth RDP/console logon screen.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1546.008",
        "T1021.001",
        "T1136.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Registry.registry_path) as registry_path values(Registry.registry_value_data) as registry_value_data values(Registry.process_name) as process_name values(Registry.user) as user from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Image File Execution Options\\\\*\" (Registry.registry_path=\"*sethc.exe*\" OR Registry.registry_path=\"*utilman.exe*\" OR Registry.registry_path=\"*osk.exe*\" OR Registry.registry_path=\"*narrator.exe*\" OR Registry.registry_path=\"*magnify.exe*\" OR Registry.registry_path=\"*displayswitch.exe*\" OR Registry.registry_path=\"*atbroker.exe*\" OR Registry.registry_path=\"*HelpPane.exe*\") Registry.registry_value_name=\"Debugger\" by Registry.dest Registry.user Registry.registry_path Registry.registry_value_name Registry.registry_value_data Registry.process_name | `drop_dm_object_name(Registry)` | where match(registry_value_data,\"(?i)(cmd\\.exe|powershell(_ise)?\\.exe|rundll32\\.exe|cscript\\.exe|wscript\\.exe|mshta\\.exe|\\\\Temp\\\\|\\\\AppData\\\\|\\\\ProgramData\\\\|\\\\Public\\\\)\") | append [| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Filesystem.file_path) as file_path values(Filesystem.process_name) as process_name from datamodel=Endpoint.Filesystem where Filesystem.action=\"modified\" (Filesystem.file_name=\"sethc.exe\" OR Filesystem.file_name=\"utilman.exe\" OR Filesystem.file_name=\"osk.exe\" OR Filesystem.file_name=\"narrator.exe\" OR Filesystem.file_name=\"magnify.exe\" OR Filesystem.file_name=\"displayswitch.exe\") Filesystem.file_path=\"*\\\\System32\\\\*\" by Filesystem.dest Filesystem.user Filesystem.file_path Filesystem.process_name | `drop_dm_object_name(Filesystem)` | where process_name!=\"TrustedInstaller.exe\" AND process_name!=\"TiWorker.exe\"] | convert ctime(firstTime) ctime(lastTime)",
      "kql": "let suspicious_debuggers = dynamic([\"cmd.exe\",\"powershell.exe\",\"powershell_ise.exe\",\"rundll32.exe\",\"cscript.exe\",\"wscript.exe\",\"mshta.exe\"]);\nlet accessibility_bins = dynamic([\"sethc.exe\",\"utilman.exe\",\"osk.exe\",\"narrator.exe\",\"magnify.exe\",\"displayswitch.exe\",\"atbroker.exe\",\"helppane.exe\"]);\nlet ifeo_hijack = DeviceRegistryEvents\n| where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has @\"\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\\"\n| where RegistryValueName =~ \"Debugger\"\n| extend hijacked = tolower(tostring(split(RegistryKey,\"\\\\\")[-1]))\n| where hijacked in (accessibility_bins)\n| where RegistryValueData has_any (suspicious_debuggers) or RegistryValueData has_any (@\"\\Temp\\\", @\"\\AppData\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\")\n| project Timestamp, DeviceName, InitiatingProcessAccountName, InitiatingProcessFileName, InitiatingProcessCommandLine, hijacked, RegistryKey, RegistryValueData, ReportId;\nlet binary_swap = DeviceFileEvents\n| where ActionType in (\"FileCreated\",\"FileModified\",\"FileRenamed\")\n| where FolderPath has @\"\\System32\\\" or FolderPath has @\"\\SysWOW64\\\"\n| where FileName in~ (accessibility_bins)\n| where InitiatingProcessFileName !in~ (\"TrustedInstaller.exe\",\"TiWorker.exe\",\"msiexec.exe\",\"poqexec.exe\")\n| project Timestamp, DeviceName, InitiatingProcessAccountName, InitiatingProcessFileName, InitiatingProcessCommandLine, FolderPath, FileName, SHA256, ReportId;\nunion ifeo_hijack, binary_swap\n| join kind=leftouter (DeviceLogonEvents | where LogonType in (\"RemoteInteractive\",\"Network\") | summarize logon_count=count(), logon_users=make_set(AccountName,25) by DeviceName, bin(Timestamp,1d)) on DeviceName\n| order by Timestamp desc",
      "rationale": "FireEye and Symantec both attribute the sticky-keys/utilman IFEO Debugger backdoor to APT3/Buckeye as a hallmark of Operation Clandestine Fox/Wolf intrusions \u2014 the group plants this so a stolen RDP credential (or even pre-auth Shift x5 at the logon screen) yields a SYSTEM cmd. Legitimate writes to IFEO Debugger keys for these specific accessibility binaries are vanishingly rare outside debugging t"
    },
    {
      "title": "[LLM] APT3 SHOTPUT/Pirpi loader: rundll32 spawned by browser/Office or scheduled task with C2 egress",
      "description": "Hunts the APT3 post-exploitation execution chain seen in Operation Clandestine Wolf and Double Tap: an IE/Flash 0-day or weaponized link spawns rundll32.exe loading a DLL (SHOTPUT / CookieCutter / Pirpi) from %TEMP%, %APPDATA%, %PROGRAMDATA% or C:\\Users\\Public\\, often persisted via schtasks, which then beacons to a SOCKS-style C2 over non-standard ports (T1090.002 / T1095).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "exploit",
      "conf": "medium",
      "techs": [
        "T1203",
        "T1218.011",
        "T1204.001",
        "T1053.005",
        "T1090.002",
        "T1095"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as process values(Processes.parent_process_name) as parent values(Processes.process_path) as process_path values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.process_name=\"rundll32.exe\" (Processes.parent_process_name IN (\"iexplore.exe\",\"outlook.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"chrome.exe\",\"firefox.exe\",\"flashplayer*.exe\",\"FlashUtil*.exe\") OR Processes.process IN (\"*\\\\AppData\\\\*\",\"*\\\\Temp\\\\*\",\"*\\\\ProgramData\\\\*\",\"*\\\\Users\\\\Public\\\\*\")) Processes.process=\"*.dll*\" by Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process Processes.process_id | `drop_dm_object_name(Processes)` | rex field=process \"(?i)(?<dll_path>[a-z]:\\\\\\\\[^\\\"\\s]+\\.dll)[,\\s]+(?<export>[A-Za-z0-9_@]+)\" | join type=inner dest process_id [| tstats summariesonly=true values(All_Traffic.dest_ip) as c2_ip values(All_Traffic.dest_port) as c2_port sum(All_Traffic.bytes_out) as bytes_out from datamodel=Network_Traffic.All_Traffic where All_Traffic.app=\"rundll32.exe\" All_Traffic.dest_category!=\"internal\" by All_Traffic.src All_Traffic.process_id | rename All_Traffic.src as dest All_Traffic.process_id as process_id] | where bytes_out > 1024 | table firstTime lastTime dest user parent_process_name process dll_path export c2_ip c2_port bytes_out | convert ctime(firstTime) ctime(lastTime)",
      "kql": "let exploit_parents = dynamic([\"iexplore.exe\",\"outlook.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"chrome.exe\",\"firefox.exe\",\"flashutil_activex.exe\",\"flashplayer.exe\",\"flashplayerplugin.exe\"]);\nlet user_writable = dynamic([@\"\\appdata\\\",@\"\\temp\\\",@\"\\programdata\\\",@\"\\users\\public\\\",@\"\\downloads\\\"]);\nlet rundll_loads =\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName =~ \"rundll32.exe\"\n| where InitiatingProcessFileName in~ (exploit_parents)\n   or ProcessCommandLine has_any (user_writable)\n| where ProcessCommandLine matches regex @\"(?i)\\.dll[,\\s\\\"']+[A-Za-z0-9_@]+\"\n| extend dll = extract(@\"(?i)([A-Za-z]:\\\\[^\\\"]+\\.dll)\", 1, ProcessCommandLine)\n| extend export = extract(@\"(?i)\\.dll[\\\",\\s]+([A-Za-z0-9_@]+)\", 1, ProcessCommandLine)\n| project Timestamp, DeviceId, DeviceName, AccountName, InitiatingProcessFileName, InitiatingProcessCommandLine, FileName, ProcessCommandLine, ProcessId, dll, export;\nlet rundll_egress =\nDeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName =~ \"rundll32.exe\"\n| where RemoteIPType == \"Public\"\n| where ActionType in (\"ConnectionSuccess\",\"ConnectionRequest\")\n| summarize bytes=sum(InitiatingProcessIntegrityLevel == \"\" and 0 or 0), c2_ips=make_set(RemoteIP,10), c2_ports=make_set(RemotePort,10), conn_count=count() by DeviceId, InitiatingProcessId, InitiatingProcessCreationTime, bin(Timestamp,1h);\nrundll_loads\n| join kind=inner (rundll_egress) on DeviceId, $left.ProcessId == $right.InitiatingProcessId\n| join kind=leftouter (\n    DeviceProcessEvents\n    | where FileName =~ \"schtasks.exe\" and ProcessCommandLine has_any (\"/create\",\"-create\")\n    | where ProcessCommandLine has_any (\"rundll32\",\".dll\")\n    | project sched_dev=DeviceId, sched_time=Timestamp, sched_cmd=ProcessCommandLine\n) on $left.DeviceId == $right.sched_dev\n| project Timestamp, DeviceName, AccountName, InitiatingProcessFileName, ProcessCommandLine, dll, export, c2_ips, c2_ports, conn_count, sched_cmd\n| order by Timestamp desc",
      "rationale": "APT3's documented exploit chains (CVE-2014-1776 IE UAF in Clandestine Fox, CVE-2015-3113 Flash UAF in Clandestine Wolf, CVE-2014-6332 in Double Tap) all converge on the same execution stub: rundll32.exe loading a SHOTPUT/Pirpi/CookieCutter DLL dropped under a user-writable directory, optionally pinned with a schtasks /create entry, then proxying C2 back out (FireEye/Recorded Future/Symantec report"
    }
  ]
}