# intel/.llm_actor_uc_cache/3b/3b7e491ea8c66c6e952705b1e9257a5228f6e401.json
{
  "name": "Agrius",
  "ucs": [
    {
      "title": "[LLM] Agrius (Pink Sandstorm / Agonizing Serpens) ASPXSpy web shell \u2014 IIS w3wp.exe spawning recon shells after public-app exploit",
      "description": "Hunts post-exploitation activity from the Agrius web shell chain: after exploiting a public-facing app (Fortinet, Exchange, Liferay, SysAid), the actor drops an ASPXSpy/MOIS variant onto the IIS server and operates via cmd.exe / PowerShell child processes spawned by w3wp.exe. This catches the very narrow handover window between T1190 \u2192 T1505.003 \u2192 T1059.003 that Agrius repeatedly uses against Israeli targets.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1190",
        "T1505.003",
        "T1059.003",
        "T1018"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdLines values(Processes.process_hash) as hashes from datamodel=Endpoint.Processes where Processes.parent_process_name=\"w3wp.exe\" AND Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"ipconfig.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"quser.exe\",\"nltest.exe\",\"arp.exe\",\"netstat.exe\",\"reg.exe\",\"sc.exe\",\"wmic.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"schtasks.exe\",\"vssadmin.exe\",\"plink.exe\",\"procdump.exe\",\"ntdsutil.exe\") by Processes.dest Processes.user Processes.parent_process Processes.process_name Processes.process | `drop_dm_object_name(Processes)` | convert ctime(firstTime) ctime(lastTime) | sort - lastTime",
      "kql": "let SuspiciousChildren = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"ipconfig.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"quser.exe\",\"nltest.exe\",\"arp.exe\",\"netstat.exe\",\"reg.exe\",\"sc.exe\",\"wmic.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"schtasks.exe\",\"vssadmin.exe\",\"plink.exe\",\"procdump.exe\",\"ntdsutil.exe\",\"7z.exe\",\"rar.exe\",\"winrar.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName =~ \"w3wp.exe\"           // IIS worker = web shell host\n| where FileName in~ (SuspiciousChildren)\n| extend WebRecon = ProcessCommandLine has_any (\"whoami\",\"net user\",\"net group\",\"ipconfig\",\"systeminfo\",\"quser\",\"nltest\",\"tasklist\")\n| extend WebExfil = ProcessCommandLine has_any (\" -hp\",\" a -p\",\"-r0\",\"DownloadString\",\"DownloadFile\",\"Invoke-WebRequest\",\"certutil -urlcache\",\"bitsadmin /transfer\")\n| project Timestamp, DeviceName, AccountName, AccountDomain,\n          AppPoolCmd = InitiatingProcessCommandLine,\n          ChildImage = FolderPath,\n          ChildCmd   = ProcessCommandLine,\n          ChildSHA256 = SHA256,\n          WebRecon, WebExfil\n| order by Timestamp desc",
      "sentinel_kql": "let SuspiciousChildren = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"ipconfig.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"quser.exe\",\"nltest.exe\",\"arp.exe\",\"netstat.exe\",\"reg.exe\",\"sc.exe\",\"wmic.exe\",\"bitsadmin.exe\",\"certutil.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"schtasks.exe\",\"vssadmin.exe\",\"plink.exe\",\"procdump.exe\",\"ntdsutil.exe\",\"7z.exe\",\"rar.exe\",\"winrar.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688                                     // process create (4688 audit must be enabled)\n| where ParentProcessName endswith \"w3wp.exe\"              // IIS worker = web shell host\n| extend ChildExe = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n| where ChildExe in (SuspiciousChildren)\n| extend WebRecon = CommandLine has_any (\"whoami\",\"net user\",\"net group\",\"ipconfig\",\"systeminfo\",\"quser\",\"nltest\",\"tasklist\")\n| extend WebExfil = CommandLine has_any (\" -hp\",\" a -p\",\"-r0\",\"DownloadString\",\"DownloadFile\",\"Invoke-WebRequest\",\"certutil -urlcache\",\"bitsadmin /transfer\")\n| project TimeGenerated, Computer, Account, SubjectUserName,\n          ParentProcessName, NewProcessName, CommandLine,\n          TokenElevationType, WebRecon, WebExfil\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "Agrius's standard playbook against Israeli/ME targets is internet-edge exploitation followed by an ASPXSpy-family web shell (Microsoft / SentinelOne / CheckPoint reporting on Agonizing Serpens 2023, Moneybird 2023, Fox Kitten/Pink Sandstorm overlap). Shell processes spawned by the w3wp.exe app pool worker are extraordinarily rare in legitimate IIS workloads (legit ASP.NET apps almost never invoke "
    },
    {
      "title": "[LLM] Agrius wiper service install \u2014 sc.exe creating service with binPath in ProgramData/Temp/Public (Apostle, DEADWOOD, MultiLayer, Moneybird, Bi",
      "description": "Hunts the service-based deployment step Agrius uses for every wiper/ransomware payload it has shipped: Apostle, DEADWOOD, Fantasy, MultiLayer, Moneybird, IPsec Helper, BiBi-Windows are all registered as Windows services with binPath pointing to a writable, non-system directory. Catches T1543.003 just before detonation.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1543.003",
        "T1059.003",
        "T1570",
        "T1036"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.process_name=\"sc.exe\" AND Processes.process=\"*create*\" AND Processes.process=\"*binpath*\" AND (Processes.process=\"*\\\\ProgramData\\\\*\" OR Processes.process=\"*\\\\Users\\\\Public\\\\*\" OR Processes.process=\"*\\\\Windows\\\\Temp\\\\*\" OR Processes.process=\"*\\\\AppData\\\\Local\\\\*\" OR Processes.process=\"*\\\\AppData\\\\Roaming\\\\*\" OR Processes.process=\"*\\\\PerfLogs\\\\*\" OR Processes.process=\"*\\\\Windows\\\\Debug\\\\*\" OR Processes.process=\"*\\\\Intel\\\\*\" OR Processes.process=\"*\\\\Recovery\\\\*\") by Processes.dest Processes.user Processes.parent_process_name Processes.parent_process Processes.process | `drop_dm_object_name(Processes)` | rex field=process \"binpath=\\\\\\\"?(?<binPath>[^\\\\\\\" ]+)\" | rex field=process \"create\\s+\\\"?(?<svcName>[^\\\\\\\" ]+)\" | convert ctime(firstTime) ctime(lastTime) | sort - lastTime",
      "kql": "let WriteableStagingPaths = dynamic([\n    @\"\\programdata\\\", @\"\\users\\public\\\", @\"\\windows\\temp\\\",\n    @\"\\appdata\\local\\\", @\"\\appdata\\roaming\\\", @\"\\perflogs\\\",\n    @\"\\windows\\debug\\\", @\"\\intel\\\", @\"\\recovery\\\", @\"\\$recycle.bin\\\"]);\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName =~ \"sc.exe\"\n| where AccountName !endswith \"$\"\n| extend cmd = tolower(ProcessCommandLine)\n| where cmd has \"create\" and cmd has \"binpath\"\n| where cmd has_any (WriteableStagingPaths)\n// Agrius binary names observed across SentinelOne / Microsoft / CheckPoint reporting\n| extend KnownAgriusName = cmd has_any (\"apostle\",\"deadwood\",\"fantasy\",\"moneybird\",\"multilayer\",\"bibi\",\"ipsechelper\",\"msupdate\",\"msupdate.exe\",\"sandals\",\"powerless\")\n| project Timestamp, DeviceName, AccountName,\n          ServiceCreateCmd = ProcessCommandLine,\n          Parent     = InitiatingProcessFileName,\n          ParentCmd  = InitiatingProcessCommandLine,\n          ParentSHA  = InitiatingProcessSHA256,\n          KnownAgriusName\n| order by Timestamp desc\n// Pivot: correlate to the dropper file write\n| join kind=leftouter (\n    DeviceFileEvents\n    | where Timestamp > ago(7d)\n    | where ActionType == \"FileCreated\"\n    | where FolderPath has_any (@\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\Windows\\Temp\\\", @\"\\PerfLogs\\\")\n    | where FileName endswith \".exe\" or FileName endswith \".dll\" or FileName endswith \".sys\"\n    | project FileTime = Timestamp, DeviceName, DroppedFile = FolderPath, DroppedSHA256 = SHA256\n  ) on DeviceName\n| where isnull(FileTime) or abs(datetime_diff('minute', Timestamp, FileTime)) <= 30",
      "sentinel_kql": "let WriteableStagingPaths = dynamic([\n    @\"\\programdata\\\", @\"\\users\\public\\\", @\"\\windows\\temp\\\",\n    @\"\\appdata\\local\\\", @\"\\appdata\\roaming\\\", @\"\\perflogs\\\",\n    @\"\\windows\\debug\\\", @\"\\intel\\\", @\"\\recovery\\\", @\"\\$recycle.bin\\\"]);\nlet ServiceCreate = SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4688                                     // requires CommandLine auditing\n    | where NewProcessName endswith \"sc.exe\"\n    | where Account !endswith \"$\"\n    | extend cmd = tolower(CommandLine)\n    | where cmd has \"create\" and cmd has \"binpath\"\n    | where cmd has_any (WriteableStagingPaths)\n    | extend KnownAgriusName = cmd has_any (\"apostle\",\"deadwood\",\"fantasy\",\"moneybird\",\"multilayer\",\"bibi\",\"ipsechelper\",\"msupdate\",\"sandals\",\"powerless\")\n    | project TimeGenerated, Computer, Account, ParentProcessName, NewProcessName, CommandLine, KnownAgriusName;\n// Companion signal: 7045 - service installed (catches Powershell New-Service / direct SCM API too)\nlet ServiceInstalled = Event\n    | where TimeGenerated > ago(7d)\n    | where EventLog == \"System\" and EventID == 7045\n    | extend EventDataLower = tolower(tostring(parse_xml(ParameterXml)))\n    | where EventDataLower has_any (WriteableStagingPaths)\n    | project TimeGenerated, Computer, RenderedDescription, ParameterXml;\nunion withsource=Source ServiceCreate, ServiceInstalled\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Agrius Wiper Service Install via sc.exe with Non-System binPath\nid: 5b3c9f81-2a4d-4e62-9a1b-7e0f8c4d6a91\nstatus: experimental\ndescription: Detects sc.exe being used to register a Windows service whose binary lives in a writeable non-system directory (ProgramData, Users\\Public, Windows\\Temp, PerfLogs). This is the standard service-installation step used by every Agrius (Pink Sandstorm / AMERICIUM / Agonizing Serpens) wiper and ransomware payload \u2014 Apostle, DEADWOOD, Fantasy, MultiLayer, Moneybird, IPsecHelper, and BiBi-Windows.\nreferences:\n  - https://attack.mitre.org/groups/G1030/\n  - https://www.sentinelone.com/labs/from-wiper-to-ransomware-the-evolution-of-agrius/\n  - https://research.checkpoint.com/2023/the-rewards-of-fantasy-iranian-apt-targets-the-diamond-industry/\n  - https://www.microsoft.com/en-us/security/blog/2023/04/18/nation-state-threat-actor-mint-sandstorm-refines-tradecraft/\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.persistence\n  - attack.privilege_escalation\n  - attack.t1543.003\n  - attack.g1030\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_img:\n    Image|endswith: '\\sc.exe'\n  selection_create:\n    CommandLine|contains|all:\n      - 'create'\n      - 'binpath'\n  selection_path:\n    CommandLine|contains:\n      - '\\ProgramData\\'\n      - '\\Users\\Public\\'\n      - '\\Windows\\Temp\\'\n      - '\\AppData\\Local\\'\n      - '\\AppData\\Roaming\\'\n      - '\\PerfLogs\\'\n      - '\\Windows\\Debug\\'\n      - '\\Intel\\'\n      - '\\Recovery\\'\n  filter_machine_account:\n    User|endswith: '$'\n  condition: selection_img and selection_create and selection_path and not filter_machine_account\nfields:\n  - Image\n  - ParentImage\n  - CommandLine\n  - User\n  - Computer\nfalsepositives:\n  - Bespoke line-of-business installers placing services under ProgramData (rare; baseline first)\n  - Antivirus / EDR self-installers under Public\\* or ProgramData\\* \u2014 allowlist by parent SHA256\nlevel: high\n",
      "rationale": "Agrius is documented (SentinelOne 'From Wiper to Ransomware' 2021, CheckPoint 'Fantasy' 2023, Microsoft 'Pink Sandstorm/MERCURY' 2022-2023) to deploy each successive wiper as a Windows service with binPath pointing to a writeable directory \u2014 Apostle was registered as 'MsUpdate', DEADWOOD as a service in %TEMP%, Fantasy via supply-chain dropper installing services into ProgramData, Moneybird and Mu"
    }
  ]
}