# intel/.llm_actor_uc_cache/4c/4c1438a24c1ddae751a835f81daffeccea0148a7.json
{
  "name": "FIN10",
  "ucs": [
    {
      "title": "[LLM] FIN10 PowerShell Empire stager + schtasks persistence chain",
      "description": "Hunts FIN10's signature post-exploitation chain: a PowerShell Empire/Meterpreter stager invoked from a non-interactive parent (wmic, services, or office) that within minutes registers a scheduled task pointing back at PowerShell or rundll32. FireEye's 2017 FIN10 report specifically called out scheduled-task persistence (often masquerading as 'Adobe' / 'Java' updaters) chained to PowerShell Empire HTTP listeners on the same host.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1059.001",
        "T1053.005",
        "T1547.001",
        "T1570"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as ps_cmd values(Processes.parent_process_name) as ps_parent from datamodel=Endpoint.Processes where Processes.process_name IN (\"powershell.exe\",\"powershell_ise.exe\") (Processes.process=\"*-NoP*\" OR Processes.process=\"*-NonI*\" OR Processes.process=\"*-w hidden*\" OR Processes.process=\"*-EncodedCommand*\" OR Processes.process=\"*FromBase64String*\" OR Processes.process=\"*IEX*\" OR Processes.process=\"*System.Net.WebClient*\" OR Processes.process=\"*DownloadString*\") by Processes.dest Processes.user Processes.parent_process_name _time span=5m | `drop_dm_object_name(Processes)` | rename dest as ps_dest user as ps_user | join type=inner ps_dest [| tstats summariesonly=true count values(Processes.process) as task_cmd from datamodel=Endpoint.Processes where Processes.process_name IN (\"schtasks.exe\",\"at.exe\") Processes.process=\"*/create*\" (Processes.process=\"*powershell*\" OR Processes.process=\"*rundll32*\" OR Processes.process=\"*regsvr32*\" OR Processes.process=\"*Adobe*\" OR Processes.process=\"*Java*\" OR Processes.process=\"*Google*Update*\" OR Processes.process=\"*\\\\Users\\\\Public\\\\*\" OR Processes.process=\"*\\\\AppData\\\\*\") by Processes.dest Processes.user _time span=5m | `drop_dm_object_name(Processes)` | rename dest as ps_dest] | where abs(_time - _time) < 600 | table firstTime ps_dest ps_user ps_parent ps_cmd task_cmd | `security_content_ctime(firstTime)`",
      "kql": "// FIN10 \u2014 PowerShell Empire stager followed by schtasks persistence on same host within 10 min\nlet WindowMin = 10m;\nlet PSStagers =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where FileName in~ (\"powershell.exe\",\"powershell_ise.exe\",\"pwsh.exe\")\n    | where AccountName !endswith \"$\"\n    | where ProcessCommandLine has_any (\"-NoP\",\"-NonI\",\"-w hidden\",\"-WindowStyle Hidden\",\"-EncodedCommand\",\"-enc \",\"FromBase64String\",\"IEX(\",\"Invoke-Expression\",\"DownloadString\",\"System.Net.WebClient\",\"Net.WebRequest\",\"Empire\")\n    | project StagerTime=Timestamp, DeviceId, DeviceName, AccountName,\n              StagerParent=InitiatingProcessFileName, StagerCmd=ProcessCommandLine;\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName in~ (\"schtasks.exe\",\"at.exe\")\n| where ProcessCommandLine has_any (\"/create\",\"-create\",\" /sc \",\" /tr \")\n| where ProcessCommandLine has_any (\"powershell\",\"rundll32\",\"regsvr32\",\"mshta\",\"wscript\",\"cscript\",\"cmd.exe\",@\"\\Users\\Public\\\",@\"\\AppData\\\",@\"\\ProgramData\\\")\n   or ProcessCommandLine matches regex @\"(?i)/tn\\s+\\\"?(adobe|java|google|microsoft|windows)?\\s*(update|updater|task|maint|defrag|run)\"\n| project TaskTime=Timestamp, DeviceId, DeviceName, AccountName,\n          TaskCmd=ProcessCommandLine, TaskParent=InitiatingProcessFileName\n| join kind=inner PSStagers on DeviceId\n| where TaskTime between (StagerTime .. StagerTime + WindowMin)\n| project StagerTime, TaskTime,\n          DelaySec = datetime_diff('second', TaskTime, StagerTime),\n          DeviceName, AccountName, StagerParent, StagerCmd, TaskParent, TaskCmd\n| order by StagerTime desc",
      "sentinel_kql": "// FIN10 \u2014 PS Empire stager + schtasks persistence (Sentinel SecurityEvent 4688)\nlet WindowMin = 10m;\nlet PSStagers =\n    SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4688\n    | where NewProcessName endswith \"powershell.exe\" or NewProcessName endswith \"powershell_ise.exe\" or NewProcessName endswith \"pwsh.exe\"\n    | where Account !endswith \"$\"\n    | where CommandLine has_any (\"-NoP\",\"-NonI\",\"-w hidden\",\"-WindowStyle Hidden\",\"-EncodedCommand\",\"-enc \",\"FromBase64String\",\"IEX(\",\"Invoke-Expression\",\"DownloadString\",\"System.Net.WebClient\",\"Empire\")\n    | project StagerTime=TimeGenerated, Computer, Account,\n              StagerParent=ParentProcessName, StagerCmd=CommandLine;\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where NewProcessName endswith \"schtasks.exe\" or NewProcessName endswith \"at.exe\"\n| where CommandLine has_any (\"/create\",\"-create\",\" /sc \",\" /tr \")\n| where CommandLine has_any (\"powershell\",\"rundll32\",\"regsvr32\",\"mshta\",\"wscript\",\"cscript\",\"cmd.exe\",@\"\\Users\\Public\\\",@\"\\AppData\\\",@\"\\ProgramData\\\")\n   or CommandLine matches regex @\"(?i)/tn\\s+\\\"?(adobe|java|google|microsoft|windows)?\\s*(update|updater|task|maint|defrag|run)\"\n| project TaskTime=TimeGenerated, Computer, Account,\n          TaskCmd=CommandLine, TaskParent=ParentProcessName\n| join kind=inner PSStagers on Computer\n| where TaskTime between (StagerTime .. StagerTime + WindowMin)\n| project StagerTime, TaskTime,\n          DelaySec = datetime_diff('second', TaskTime, StagerTime),\n          Computer, Account, StagerParent, StagerCmd, TaskParent, TaskCmd\n| order by StagerTime desc",
      "sigma_yaml": "",
      "rationale": "FIN10's installation tradecraft is documented (FireEye, June 2017) as PowerShell Empire stagers paired with schtasks-based persistence using lure names like 'Adobe Update' or 'Java Update' pointing at PowerShell/rundll32 payloads in C:\\Users\\Public or %AppData%. Most APTs use one or the other; the same-host, sub-10-minute pairing of PS Empire-style cmdline plus schtasks /create with a PS/rundll32 "
    },
    {
      "title": "[LLM] FIN10 RDP lateral movement to data-staging share with valid stolen accounts",
      "description": "Hunts FIN10's post-foothold lateral movement: an interactive RDP logon using a non-machine account from an internal hop-host followed within an hour by archive/staging activity (.7z/.rar/.zip creation) on the destination \u2014 the FireEye-documented FIN10 path from initial compromise to bulk data theft used to extort gaming/mining victims. Catches the actor's preference for valid accounts + RDP rather than C2-driven lateral movement.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "medium",
      "techs": [
        "T1021.001",
        "T1078",
        "T1078.003",
        "T1033",
        "T1070.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as logonTime from datamodel=Authentication where Authentication.action=success Authentication.authentication_method=\"Kerberos\" OR Authentication.authentication_method=\"NTLM\" Authentication.signature_id IN (4624) Authentication.logon_type=10 NOT Authentication.user IN (\"*$\",\"SYSTEM\",\"ANONYMOUS LOGON\") by Authentication.dest Authentication.user Authentication.src _time span=1h | `drop_dm_object_name(Authentication)` | rename dest as host user as actor | join type=inner host actor [| tstats summariesonly=true count values(Filesystem.file_name) as staged_files values(Filesystem.file_path) as staged_paths from datamodel=Endpoint.Filesystem where Filesystem.action=created (Filesystem.file_name=\"*.7z\" OR Filesystem.file_name=\"*.rar\" OR Filesystem.file_name=\"*.zip\" OR Filesystem.file_name=\"*.cab\") (Filesystem.file_path=\"*\\\\Users\\\\Public\\\\*\" OR Filesystem.file_path=\"*\\\\AppData\\\\*\" OR Filesystem.file_path=\"*\\\\ProgramData\\\\*\" OR Filesystem.file_path=\"*\\\\Temp\\\\*\" OR Filesystem.file_path=\"*\\\\Windows\\\\Tasks\\\\*\") by Filesystem.dest Filesystem.user _time span=1h | `drop_dm_object_name(Filesystem)` | rename dest as host user as actor] | where abs(_time - logonTime) < 3600 AND _time >= logonTime | table logonTime _time host actor src staged_files staged_paths count | `security_content_ctime(logonTime)`",
      "kql": "// FIN10 \u2014 RDP logon (LogonType 10) followed by archive staging on same host within 60m\nlet StageWindow = 60m;\nlet RdpLogons =\n    DeviceLogonEvents\n    | where Timestamp > ago(7d)\n    | where ActionType == \"LogonSuccess\"\n    | where LogonType == \"RemoteInteractive\"\n    | where AccountName !endswith \"$\"\n    | where AccountName !in~ (\"system\",\"local service\",\"network service\",\"anonymous logon\")\n    | where RemoteIPType in (\"Private\",\"Loopback\")    // internal hop, not external RDP gateway\n    | project LogonTime=Timestamp, DeviceId, DeviceName, AccountName,\n              SourceHost=RemoteDeviceName, SourceIP=RemoteIP, LogonId;\nDeviceFileEvents\n| where Timestamp > ago(7d)\n| where ActionType == \"FileCreated\"\n| where FileName endswith \".7z\" or FileName endswith \".rar\" or FileName endswith \".zip\" or FileName endswith \".cab\" or FileName matches regex @\"(?i)\\.(00[0-9]|r0[0-9])$\"\n| where FolderPath has_any (@\"\\Users\\Public\\\", @\"\\AppData\\\", @\"\\ProgramData\\\", @\"\\Temp\\\", @\"\\Windows\\Tasks\\\", @\"\\PerfLogs\\\")\n| where InitiatingProcessFileName in~ (\"7z.exe\",\"7za.exe\",\"rar.exe\",\"winrar.exe\",\"makecab.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\")\n| project StageTime=Timestamp, DeviceId, DeviceName,\n          ArchiveFile=FileName, ArchivePath=FolderPath,\n          StageActor=InitiatingProcessAccountName,\n          StageCmd=InitiatingProcessCommandLine\n| join kind=inner RdpLogons on DeviceId\n| where StageTime between (LogonTime .. LogonTime + StageWindow)\n| where StageActor =~ AccountName\n| project LogonTime, StageTime,\n          MinutesAfterLogon = datetime_diff('minute', StageTime, LogonTime),\n          DeviceName, AccountName, SourceHost, SourceIP,\n          ArchiveFile, ArchivePath, StageCmd\n| order by LogonTime desc",
      "sentinel_kql": "// FIN10 \u2014 RDP logon (4624 LogonType 10) + archive staging on same host within 60m\nlet StageWindow = 60m;\nlet RdpLogons =\n    SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4624\n    | where LogonType == 10\n    | where TargetUserName !endswith \"$\"\n    | where TargetUserName !in~ (\"SYSTEM\",\"ANONYMOUS LOGON\",\"LOCAL SERVICE\",\"NETWORK SERVICE\")\n    | where ipv4_is_private(IpAddress)\n    | project LogonTime=TimeGenerated, Computer, Account=TargetUserName,\n              SourceIP=IpAddress, SourceHost=WorkstationName, TargetLogonId;\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where (CommandLine has_any (\"7z.exe\",\"7za.exe\",\"rar.exe\",\"winrar.exe\",\"makecab.exe\")\n         or NewProcessName endswith \"7z.exe\" or NewProcessName endswith \"rar.exe\" or NewProcessName endswith \"winrar.exe\")\n| where CommandLine has_any (\" a \",\" -m\",\" -hp\",\" -p\",\"-r \")\n| where CommandLine has_any (\".7z\",\".rar\",\".zip\",\".cab\",@\"\\Users\\Public\\\",@\"\\AppData\\\",@\"\\ProgramData\\\",@\"\\Temp\\\",@\"\\PerfLogs\\\")\n| project StageTime=TimeGenerated, Computer, Account=SubjectUserName,\n          StageCmd=CommandLine, StageProc=NewProcessName\n| join kind=inner RdpLogons on Computer\n| where StageTime between (LogonTime .. LogonTime + StageWindow)\n| where Account =~ tostring(split(RdpLogons.Account,\"@\")[0])\n| project LogonTime, StageTime,\n          MinutesAfterLogon = datetime_diff('minute', StageTime, LogonTime),\n          Computer, Account, SourceHost, SourceIP,\n          StageProc, StageCmd\n| order by LogonTime desc",
      "sigma_yaml": "",
      "rationale": "FIN10's intrusion outcome \u2014 per the FireEye 2017 cluster report \u2014 is bulk data theft used to extort North American mining and casino victims. Their lateral movement pattern is RDP with stolen valid accounts (T1078 + T1021.001), not C2-driven SMB pivots, and they stage exfil archives in writable directories like C:\\Users\\Public or C:\\PerfLogs before exfiltration. The chained pattern (interactive RD"
    }
  ]
}