# intel/.llm_actor_uc_cache/82/82fa9354f47d7d39dac2117a3b17f07da793ae90.json
{
  "name": "Group5",
  "ucs": [
    {
      "title": "[LLM] Group5 commodity RAT install chain \u2014 script/document parent drops EXE to AppData and self-registers Run key (njRAT / NanoCore)",
      "description": "Detects the canonical install footprint of njRAT and NanoCore \u2014 the two commodity RATs Citizen Lab tied to Group5's Syrian-opposition spearphishing campaigns. The detection ties together the spearphish-spawned dropper (script host, Office app, or installer) creating an .exe in AppData with the same binary being registered as HKCU Run/RunOnce persistence within minutes \u2014 the deterministic install pattern of both njRAT (HKCU\\...\\Run, AppData\\Roaming) and NanoCore (HKCU\\...\\Run, AppData\\Roaming\\<GUID>\\<random>.exe).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1547.001",
        "T1204.002",
        "T1027.013",
        "T1056.001",
        "T1113"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true min(_time) as drop_time values(Filesystem.file_path) as dropped_file values(Filesystem.file_name) as dropped_name values(Filesystem.process_name) as dropper_proc values(Filesystem.user) as user from datamodel=Endpoint.Filesystem where Filesystem.action=created AND (Filesystem.file_path=\"*\\\\AppData\\\\Roaming\\\\*\" OR Filesystem.file_path=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Filesystem.file_path=\"*\\\\AppData\\\\Local\\\\*\") AND Filesystem.file_name=\"*.exe\" AND Filesystem.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"hh.exe\",\"rundll32.exe\") by Filesystem.dest | rename Filesystem.dest AS dest | join type=inner dest [| tstats summariesonly=true min(_time) as reg_time values(Registry.registry_path) as reg_key values(Registry.registry_value_name) as reg_value_name values(Registry.registry_value_data) as reg_value_data values(Registry.process_name) as reg_writer from datamodel=Endpoint.Registry where Registry.action=modified AND (Registry.registry_path=\"*\\\\CurrentVersion\\\\Run*\" OR Registry.registry_path=\"*\\\\CurrentVersion\\\\RunOnce*\") AND (Registry.registry_value_data=\"*\\\\AppData\\\\Roaming\\\\*\" OR Registry.registry_value_data=\"*\\\\AppData\\\\Local\\\\*\") AND Registry.registry_value_data=\"*.exe\" by Registry.dest | rename Registry.dest AS dest] | eval delay_sec = reg_time - drop_time | where delay_sec >= 0 AND delay_sec <= 600 | eval persisted_path_matches_drop = if(like(reg_value_data, \"%\".dropped_name.\"%\"), 1, 0) | where persisted_path_matches_drop=1 | table drop_time, reg_time, delay_sec, dest, user, dropper_proc, dropped_file, reg_key, reg_value_name, reg_value_data, reg_writer | sort - drop_time",
      "kql": "// Group5 \u2014 njRAT / NanoCore install chain: AppData drop + same-binary Run-key persistence within 10 minutes\nlet lookback = 7d;\nlet WindowMinutes = 10m;\nlet SuspiciousParents = dynamic([\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"hh.exe\",\"rundll32.exe\",\"explorer.exe\"]);\nlet DroppedExe = DeviceFileEvents\n    | where Timestamp > ago(lookback)\n    | where ActionType == \"FileCreated\"\n    | where FolderPath has_any (@\"\\AppData\\Roaming\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\AppData\\Local\\\")\n    | where FileName endswith \".exe\"\n    | where InitiatingProcessFileName in~ (SuspiciousParents)\n    | where InitiatingProcessAccountName !endswith \"$\"          // exclude machine accounts\n    | project DropTime = Timestamp, DeviceId, DeviceName,\n              DroppedFolder = FolderPath, DroppedName = FileName,\n              DroppedSHA256 = SHA256,\n              DropperProc = InitiatingProcessFileName,\n              DropperCmd  = InitiatingProcessCommandLine,\n              DropperUser = InitiatingProcessAccountName;\nlet RunKeyPersist = DeviceRegistryEvents\n    | where Timestamp > ago(lookback)\n    | where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n    | where RegistryKey has_any (@\"\\Microsoft\\Windows\\CurrentVersion\\Run\",\n                                 @\"\\Microsoft\\Windows\\CurrentVersion\\RunOnce\")\n    | where RegistryValueData has_any (@\"\\AppData\\Roaming\\\", @\"\\AppData\\Local\\\")\n        and RegistryValueData endswith \".exe\"\n    | project RegTime = Timestamp, DeviceId,\n              RegistryKey, RegistryValueName, RegistryValueData,\n              RegInitiator    = InitiatingProcessFileName,\n              RegInitiatorCmd = InitiatingProcessCommandLine;\nDroppedExe\n| join kind=inner RunKeyPersist on DeviceId\n| where RegTime between (DropTime .. DropTime + WindowMinutes)\n| where RegistryValueData has DroppedName                       // persisted value points at the just-dropped binary\n| project DropTime, RegTime,\n          DelaySec = datetime_diff('second', RegTime, DropTime),\n          DeviceName, DropperUser,\n          DropperProc, DropperCmd,\n          DroppedFolder, DroppedName, DroppedSHA256,\n          RegistryKey, RegistryValueName, RegistryValueData,\n          RegInitiator, RegInitiatorCmd\n| order by DropTime desc",
      "sentinel_kql": "// Group5 \u2014 njRAT / NanoCore install chain (Sentinel ASIM port; assumes MDE/Sysmon ingestion via ImFileEvent + ImRegistryEvent)\nlet lookback = 7d;\nlet WindowMinutes = 10m;\nlet SuspiciousParents = dynamic([\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"hh.exe\",\"rundll32.exe\",\"explorer.exe\"]);\nlet DroppedExe = ImFileEvent\n    | where TimeGenerated > ago(lookback)\n    | where EventType == \"FileCreated\"\n    | where TargetFilePath has_any (@\"\\AppData\\Roaming\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\AppData\\Local\\\")\n    | where TargetFileName endswith \".exe\"\n    | where ActingProcessName has_any (SuspiciousParents)\n    | where ActorUsername !endswith \"$\"\n    | project DropTime = TimeGenerated, DvcHostname,\n              DroppedFilePath = TargetFilePath,\n              DroppedName = TargetFileName,\n              DroppedSHA256 = TargetFileSHA256,\n              DropperProc = ActingProcessName,\n              DropperCmd  = ActingProcessCommandLine,\n              DropperUser = ActorUsername;\nlet RunKeyPersist = ImRegistryEvent\n    | where TimeGenerated > ago(lookback)\n    | where EventType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n    | where RegistryKey has_any (@\"\\Microsoft\\Windows\\CurrentVersion\\Run\",\n                                 @\"\\Microsoft\\Windows\\CurrentVersion\\RunOnce\")\n    | where RegistryValueData has_any (@\"\\AppData\\Roaming\\\", @\"\\AppData\\Local\\\")\n        and RegistryValueData endswith \".exe\"\n    | project RegTime = TimeGenerated, DvcHostname,\n              RegistryKey, RegistryValue, RegistryValueData,\n              RegInitiator    = ActingProcessName,\n              RegInitiatorCmd = ActingProcessCommandLine;\nDroppedExe\n| join kind=inner RunKeyPersist on DvcHostname\n| where RegTime between (DropTime .. DropTime + WindowMinutes)\n| where RegistryValueData has DroppedName\n| project DropTime, RegTime,\n          DelaySec = datetime_diff('second', RegTime, DropTime),\n          DvcHostname, DropperUser,\n          DropperProc, DropperCmd,\n          DroppedFilePath, DroppedName, DroppedSHA256,\n          RegistryKey, RegistryValue, RegistryValueData,\n          RegInitiator, RegInitiatorCmd\n| order by DropTime desc",
      "sigma_yaml": "",
      "rationale": "Group5 (G0043) is a low-sophistication, suspected-Iranian-nexus actor whose entire intrusion toolkit is commodity RATs \u2014 njRAT, NanoCore and Android DroidJack \u2014 distributed via spearphishing and watering-hole drops to Syrian-opposition targets. Citizen Lab's reverse engineering of the Group5 samples documented two install footprints that this UC anchors on: njRAT writing itself to %AppData%\\Roamin"
    }
  ]
}