# rule_packs/sentinel/UC_SERVICE_PERSIST.json
{
  "schema": "https://schema.management.azure.com/schemas/2019-04-01/deploymentTemplate.json#",
  "contentVersion": "1.0.0.0",
  "resources": [
    {
      "type": "Microsoft.OperationalInsights/workspaces/providers/alertRules",
      "apiVersion": "2022-12-01-preview",
      "name": "[concat(parameters('workspaceName'), '/Microsoft.SecurityInsights/UC_SERVICE_PERSIST')]",
      "kind": "Scheduled",
      "properties": {
        "displayName": "Service install for persistence \u2014 sc.exe / new service registry write",
        "description": "Service install with binPath pointing to user-writeable path or LOLBin. Detects via process telemetry (sc.exe create) and registry (HKLM\\SYSTEM\\CurrentControlSet\\Services\\* writes).",
        "severity": "Low",
        "enabled": false,
        "query": "// EventID 4697 is \"Service installed\" \u2014 direct signal, no cmdline-regex needed.\nunion\n    ( SecurityEvent\n        | where TimeGenerated > ago(7d)\n        | where EventID == 4697\n        | where SubjectUserName !endswith \"$\"\n        | project TimeGenerated, Computer, Account = SubjectUserName,\n                  Detail = strcat(\"ServiceName=\", ServiceName, \" | Image=\", ServiceFileName) ),\n    ( SecurityEvent\n        | where TimeGenerated > ago(7d)\n        | where EventID == 4688\n        | where Account !endswith \"$\"\n        | where NewProcessName endswith \"\\sc.exe\" and CommandLine has \"create\"\n        | where CommandLine matches regex @\"(?i)(\\Users\\|\\AppData\\|\\ProgramData\\|\\Temp\\)\"\n        | project TimeGenerated, Computer, Account,\n                  Detail = strcat(\"4688 sc.exe create: \", CommandLine) )\n| order by TimeGenerated desc\n",
        "queryFrequency": "PT1D",
        "queryPeriod": "PT24H",
        "triggerOperator": "GreaterThan",
        "triggerThreshold": 0,
        "tactics": [],
        "techniques": [
          "T1543.003"
        ],
        "alertRuleTemplateName": "UC_SERVICE_PERSIST",
        "customDetails": {
          "tier": "hunting",
          "fp_rate_estimate": "high",
          "source_url": "https://github.com/Virtualhaggis/usecaseintel/blob/main/use_cases/install/UC_SERVICE_PERSIST.yml"
        }
      }
    }
  ],
  "parameters": {
    "workspaceName": {
      "type": "string"
    }
  }
}