# intel/.llm_actor_uc_cache/14/14521069ca470a79234dd7c0a96662e330c18dfb.json
{
  "name": "Leviathan",
  "ucs": [
    {
      "title": "[LLM] APT40 (Leviathan) spearphish chain: maldoc \u2192 WMI/PowerShell exec \u2192 suspicious BITS staging",
      "description": "Hunts Leviathan/APT40's signature initial access pattern: a maritime/defense-themed spearphishing attachment (T1566.001) executed by the user (T1204.002), spawning Office \u2192 WMI/PowerShell child processes (T1059.001, T1047) followed by BITS-based payload retrieval (T1197) and DLL-side-loaded second-stage (T1574.002 / T1218). This combination of Office \u2192 wmic/powershell \u2192 bitsadmin within a short window is highly characteristic of MUDCARP/TEMP.Periscope intrusions documented by CISA AA21-200A and Mandiant.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "exploit",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.001",
        "T1047",
        "T1197",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as processes values(Processes.parent_process) as parents values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\") Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\",\"wmic.exe\",\"cmd.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\") by Processes.dest Processes.user Processes.parent_process_name Processes.process_name _time span=1h\n| `drop_dm_object_name(Processes)`\n| join type=inner dest [\n    | tstats summariesonly=t count from datamodel=Endpoint.Processes where Processes.process_name=\"bitsadmin.exe\" (Processes.process=\"*transfer*\" OR Processes.process=\"*addfile*\" OR Processes.process=\"*http*\") by Processes.dest Processes.parent_process_name Processes.process _time span=1h\n    | `drop_dm_object_name(Processes)`\n    | rename count as bits_count process as bits_cmd parent_process_name as bits_parent\n  ]\n| where bits_count >= 1\n| eval risk_score = count + (bits_count*5)\n| where risk_score >= 6\n| `security_content_ctime(firstTime)` \n| table firstTime dest user parent_process_name process_name processes bits_parent bits_cmd risk_score",
      "kql": "let officeChildren = dynamic([\"powershell.exe\",\"pwsh.exe\",\"wmic.exe\",\"cmd.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"cscript.exe\",\"wscript.exe\"]);\nlet officeParents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\"]);\nlet stage1 = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (officeParents)\n| where FileName in~ (officeChildren)\n| project s1Time=Timestamp, DeviceId, DeviceName, AccountName, s1Parent=InitiatingProcessFileName, s1Child=FileName, s1Cmd=ProcessCommandLine;\nlet stage2 = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName =~ \"bitsadmin.exe\" or (FileName =~ \"powershell.exe\" and ProcessCommandLine has_any (\"Start-BitsTransfer\",\"BITS\"))\n| where ProcessCommandLine has_any (\"http://\",\"https://\",\"transfer\",\"addfile\")\n| project s2Time=Timestamp, DeviceId, s2Cmd=ProcessCommandLine, s2Parent=InitiatingProcessFileName;\nstage1\n| join kind=inner stage2 on DeviceId\n| where s2Time between (s1Time .. s1Time + 1h)\n| extend chain = strcat(s1Parent, \" -> \", s1Child, \" -> bits\")\n| summarize firstSeen=min(s1Time), lastSeen=max(s2Time), commands=make_set(pack(\"office\",s1Cmd,\"bits\",s2Cmd),10) by DeviceName, AccountName, chain\n| order by firstSeen desc",
      "rationale": "CISA AA21-200A and Mandiant reporting on TEMP.Periscope/APT40 explicitly call out the Office-macro \u2192 WMI/PowerShell \u2192 BITS staging chain (used to drop AIRBREAK/MURKYTOP/HOMEFRY against maritime, defense, and academic targets). The narrow 1-hour pivot window plus the BITS leg differentiates this from generic phishing-detection logic."
    },
    {
      "title": "[LLM] APT40 web-shell + scheduled-task persistence on Internet-facing servers (CHINACHOPPER/PIPEMON pattern)",
      "description": "Hunts Leviathan's signature post-exploitation pattern after Exploit-Public-Facing-Application (T1190) intrusions on Exchange/Confluence/Citrix: w3wp.exe / java / tomcat spawning cmd/powershell (web shell behaviour, T1505.003) followed shortly by scheduled-task creation (T1053.005) and credential dumping via comsvcs/rundll32 minidump (T1003.001). APT40 has been repeatedly observed using CHINACHOPPER-style ASPX shells alongside this exact lateral pivot.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1190",
        "T1505.003",
        "T1059.001",
        "T1003.001",
        "T1053.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdlines values(Processes.process_name) as child_procs from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"w3wp.exe\",\"httpd.exe\",\"tomcat*.exe\",\"java.exe\",\"nginx.exe\") Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"net.exe\",\"net1.exe\",\"whoami.exe\",\"systeminfo.exe\",\"ipconfig.exe\",\"schtasks.exe\",\"rundll32.exe\") by Processes.dest Processes.parent_process_name Processes.user _time span=30m\n| `drop_dm_object_name(Processes)`\n| eval webshell_score = mvcount(split(cmdlines,\"\\n\"))\n| join type=inner dest [\n    | tstats summariesonly=t count from datamodel=Endpoint.Processes where (Processes.process_name=\"schtasks.exe\" Processes.process=\"*/create*\") OR (Processes.process_name=\"rundll32.exe\" Processes.process=\"*comsvcs*MiniDump*\") OR (Processes.process_name=\"powershell.exe\" Processes.process=\"*lsass*\") by Processes.dest Processes.process Processes.process_name _time span=30m\n    | `drop_dm_object_name(Processes)`\n    | stats values(process) as followon_cmds count as followon_count by dest\n  ]\n| where followon_count >= 1\n| table firstTime dest user parent_process_name child_procs cmdlines followon_cmds followon_count\n| sort - firstTime",
      "kql": "let webProcs = dynamic([\"w3wp.exe\",\"httpd.exe\",\"tomcat.exe\",\"tomcat9.exe\",\"java.exe\",\"nginx.exe\"]);\nlet sus = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"net.exe\",\"net1.exe\",\"whoami.exe\",\"systeminfo.exe\",\"ipconfig.exe\",\"schtasks.exe\",\"rundll32.exe\",\"certutil.exe\"]);\nlet webshell = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (webProcs)\n| where FileName in~ (sus)\n| project wsTime=Timestamp, DeviceId, DeviceName, wsParent=InitiatingProcessFileName, wsChild=FileName, wsCmd=ProcessCommandLine, AccountName;\nlet followon = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where (FileName =~ \"schtasks.exe\" and ProcessCommandLine has \"/create\")\n   or (FileName =~ \"rundll32.exe\" and ProcessCommandLine has_all (\"comsvcs\",\"MiniDump\"))\n   or (FileName =~ \"powershell.exe\" and ProcessCommandLine has_any (\"lsass\",\"Out-Minidump\",\"Invoke-Mimikatz\"))\n   or (FileName =~ \"ntdsutil.exe\")\n| project foTime=Timestamp, DeviceId, foCmd=ProcessCommandLine, foFile=FileName;\nwebshell\n| join kind=inner followon on DeviceId\n| where foTime between (wsTime .. wsTime + 30m)\n| extend chain = strcat(wsParent, \" -> \", wsChild, \" -> \", foFile)\n| summarize firstSeen=min(wsTime), lastSeen=max(foTime), webshell_cmds=make_set(wsCmd,10), followon_cmds=make_set(foCmd,10) by DeviceName, AccountName, chain\n| join kind=leftouter (\n    DeviceFileEvents\n    | where Timestamp > ago(7d)\n    | where FolderPath has_any (\"\\\\inetpub\\\\wwwroot\\\\\",\"\\\\webapps\\\\\",\"\\\\Exchange Server\\\\\",\"\\\\OWA\\\\\",\"\\\\Confluence\\\\\")\n    | where FileName endswith \".aspx\" or FileName endswith \".jsp\" or FileName endswith \".ashx\" or FileName endswith \".asp\"\n    | summarize webshell_files=make_set(FolderPath,10) by DeviceName\n  ) on DeviceName\n| order by firstSeen desc",
      "rationale": "APT40/Leviathan (BRONZE MOHAWK, Kryptonite Panda) is documented by CISA, ACSC, NCSC and CrowdStrike for exploiting public-facing apps (Exchange ProxyLogon/ProxyShell, Atlassian Confluence, Citrix) to drop CHINACHOPPER-family ASPX/JSP shells, then immediately staging schtasks persistence and comsvcs/LSASS dumping for credential access. The w3wp/tomcat\u2192shell + 30-min schtasks/lsass pivot is the acto"
    }
  ]
}