# intel/.llm_actor_uc_cache/97/97c1ac0bb2759949156d26456b2ca46b5301517c.json
{
  "name": "Moonstone Sleet",
  "ucs": [
    {
      "title": "[LLM] Moonstone Sleet (Storm-1789) trojanized npm package post-install code execution",
      "description": "Hunts the developer-targeted supply-chain chain Moonstone Sleet uses to compromise software/crypto engineers via malicious npm packages (e.g. harthat-api, harthat-hash, scratchcurrencytoken, node-hide-console-win), where node/npm spawns native LOLBins or PowerShell during install to stage second-stage payloads. This is a near-signature pattern for Storm-1789's developer-recruitment lure as documented by Microsoft Threat Intelligence.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1195.002",
        "T1204.002",
        "T1059.007",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.process_path) as path from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"node.exe\",\"npm.cmd\",\"npm.exe\",\"npx.cmd\",\"yarn.cmd\",\"yarn.exe\") AND Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"curl.exe\",\"certutil.exe\",\"schtasks.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"bitsadmin.exe\") by host Processes.user Processes.parent_process Processes.parent_process_name Processes.process_name Processes.process Processes.process_path | `drop_dm_object_name(Processes)` | eval suspicious=if(match(process,\"(?i)(invoke-webrequest|downloadstring|-enc |frombase64|iex |webclient|net\\.webclient|http(s)?:\\/\\/|appdata|programdata|temp)\"),1,0), known_pkg=if(match(parent_process,\"(?i)(harthat-api|harthat-hash|node-hide-console-win|scratchcurrencytoken|postinstall|preinstall)\"),1,0) | where suspicious=1 OR known_pkg=1 | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` | sort - lastTime",
      "kql": "DeviceProcessEvents\n| where Timestamp > ago(30d)\n| where InitiatingProcessFileName in~ (\"node.exe\",\"npm.cmd\",\"npm.exe\",\"npx.cmd\",\"yarn.cmd\",\"yarn.exe\")\n| where FileName in~ (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"curl.exe\",\"certutil.exe\",\"schtasks.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"bitsadmin.exe\")\n| extend SuspChild = ProcessCommandLine has_any (\"Invoke-WebRequest\",\"DownloadString\",\"FromBase64\",\"-enc \",\"-EncodedCommand\",\"IEX \",\"WebClient\",\"http://\",\"https://\",\"\\\\AppData\\\\\",\"\\\\ProgramData\\\\\",\"\\\\Temp\\\\\")\n| extend KnownPkg = InitiatingProcessCommandLine has_any (\"harthat-api\",\"harthat-hash\",\"node-hide-console-win\",\"scratchcurrencytoken\",\"postinstall\",\"preinstall\",\"npm install\",\"npm i \")\n| where SuspChild or KnownPkg\n| join kind=leftouter (\n    DeviceNetworkEvents\n    | where Timestamp > ago(30d)\n    | where InitiatingProcessFileName in~ (\"node.exe\",\"powershell.exe\",\"curl.exe\",\"certutil.exe\")\n    | where RemoteUrl has_any (\"detankzone\",\"starglowventures\",\"ccwaterfall\",\"blocknovas\",\"matuchok\")\n      or RemoteIPType == \"Public\"\n    | project DeviceId, RemoteUrl, RemoteIP, NetTime=Timestamp\n) on DeviceId\n| project Timestamp, DeviceName, AccountName, InitiatingProcessFileName, InitiatingProcessCommandLine, FileName, ProcessCommandLine, FolderPath, RemoteUrl, RemoteIP\n| order by Timestamp desc",
      "rationale": "Storm-1789 / Moonstone Sleet is one of the very few state actors (alongside Jade Sleet/TraderTraitor) using malicious npm packages to target Web3 and software developers. Legitimate npm installs almost never spawn powershell.exe with -enc, certutil, or schtasks from a postinstall hook. Anchoring on the node/npm parent + suspicious LOLBin child + Storm-1789 package names yields very high fidelity v"
    },
    {
      "title": "[LLM] Moonstone Sleet trojanized PuTTY / DeTankZone game loader chain (SplitLoader / YouieLoader)",
      "description": "Detects the actor's signature lure where a fake-recruiter persona ships a tar.lz / ZIP containing a trojanized PuTTY (IsoBurst) or the DeTankZone / DeTankWar / DeFiTankWar game; once executed from a user-writable path the binary side-loads SplitLoader/YouieLoader, then establishes scheduled-task or Run-key persistence and beacons over HTTPS. The combination of game/PuTTY-named binary launched from Downloads/Temp + immediate persistence is uniquely Moonstone Sleet.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.003",
        "T1204.002",
        "T1547.001",
        "T1053.005",
        "T1027.013",
        "T1071.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdlines values(Processes.process_path) as paths values(Processes.process_name) as procs from datamodel=Endpoint.Processes where (match(Processes.process_path,\"(?i)\\\\\\\\Users\\\\\\\\[^\\\\\\\\]+\\\\\\\\(Downloads|Desktop|AppData\\\\\\\\Local\\\\\\\\Temp|Public)\\\\\\\\\") AND Processes.process_name IN (\"putty.exe\",\"isoburst.exe\",\"detankzone.exe\",\"detankwar.exe\",\"defitankwar.exe\",\"deltank.exe\",\"setup.exe\",\"installer.exe\")) OR (Processes.parent_process_name IN (\"putty.exe\",\"isoburst.exe\",\"detankzone.exe\",\"detankwar.exe\",\"defitankwar.exe\") AND Processes.process_name IN (\"rundll32.exe\",\"regsvr32.exe\",\"schtasks.exe\",\"reg.exe\",\"powershell.exe\",\"cmd.exe\",\"wscript.exe\",\"mshta.exe\")) by host Processes.user Processes.parent_process_name Processes.parent_process Processes.process_name Processes.process Processes.process_path | `drop_dm_object_name(Processes)` | eventstats dc(process_name) as procVariety values(parent_process_name) as parents by host user | where procVariety>=2 OR like(parent_process_name,\"%putty%\") OR like(parent_process_name,\"%tank%\") | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` | sort - lastTime",
      "kql": "let LureBins = dynamic([\"putty.exe\",\"isoburst.exe\",\"detankzone.exe\",\"detankwar.exe\",\"defitankwar.exe\",\"deltank.exe\"]);\nlet UserPaths = dynamic([@\"\\Users\\\",@\"\\Downloads\\\",@\"\\Desktop\\\",@\"\\AppData\\Local\\Temp\\\",@\"\\Users\\Public\\\"]);\nlet Stage1 =\n    DeviceProcessEvents\n    | where Timestamp > ago(30d)\n    | where FileName in~ (LureBins) and FolderPath has_any (UserPaths)\n    | project T1=Timestamp, DeviceId, DeviceName, AccountName, Stage1Proc=FileName, Stage1Path=FolderPath, Stage1Cmd=ProcessCommandLine, Stage1PID=ProcessId;\nlet Stage2 =\n    DeviceProcessEvents\n    | where Timestamp > ago(30d)\n    | where InitiatingProcessFileName in~ (LureBins)\n       and FileName in~ (\"rundll32.exe\",\"regsvr32.exe\",\"schtasks.exe\",\"reg.exe\",\"powershell.exe\",\"cmd.exe\",\"wscript.exe\",\"mshta.exe\")\n    | project T2=Timestamp, DeviceId, Stage2Proc=FileName, Stage2Cmd=ProcessCommandLine, InitProc=InitiatingProcessFileName;\nlet Stage3Persist =\n    union\n    (DeviceRegistryEvents\n     | where Timestamp > ago(30d)\n     | where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n     | where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\")\n     | where InitiatingProcessFileName in~ (LureBins) or InitiatingProcessParentFileName in~ (LureBins)\n     | project T3=Timestamp, DeviceId, PersistType=\"RunKey\", PersistDetail=strcat(RegistryKey,\"\\\\\",RegistryValueName,\" => \",RegistryValueData)),\n    (DeviceProcessEvents\n     | where Timestamp > ago(30d)\n     | where FileName =~ \"schtasks.exe\" and ProcessCommandLine has_any (\"/create\",\"/CREATE\")\n     | where InitiatingProcessFileName in~ (LureBins) or InitiatingProcessParentFileName in~ (LureBins)\n     | project T3=Timestamp, DeviceId, PersistType=\"SchedTask\", PersistDetail=ProcessCommandLine);\nStage1\n| join kind=inner (Stage2) on DeviceId\n| where T2 between (T1 .. T1 + 30m)\n| join kind=leftouter (Stage3Persist) on DeviceId\n| where isnull(T3) or T3 between (T1 .. T1 + 2h)\n| project T1, DeviceName, AccountName, Stage1Proc, Stage1Path, Stage1Cmd, Stage2Proc, Stage2Cmd, PersistType, PersistDetail\n| order by T1 desc",
      "rationale": "Microsoft and Mandiant both attribute the DeTankZone/DeTankWar fake-game and the trojanized PuTTY (IsoBurst) campaigns specifically to Moonstone Sleet / Storm-1789 \u2014 these are not shared with broader Lazarus clusters since the 2023 split. Real PuTTY almost never executes from Downloads/Temp and then spawns rundll32/schtasks, and a binary literally named detankzone.exe spawning persistence is essen"
    }
  ]
}