# intel/.llm_actor_uc_cache/47/4785f7d7cedc7f19071eec93b2af071e2c75faf4.json
{
  "name": "Embargo",
  "ucs": [
    {
      "title": "[LLM] Embargo MDeployer + MS4Killer chain: Safe Mode reboot for BYOVD EDR neutralization",
      "description": "Hunts Embargo's signature pre-encryption sequence where the Rust-based MDeployer loader stages MS4Killer, modifies the boot configuration to force a Safe Mode (with networking) reboot, and re-launches itself as a service to terminate AV/EDR via a vulnerable signed driver. ESET's analysis of Embargo (suspected ALPHV/BlackCat successor) documents this exact bcdedit\u2192shutdown\u2192safe-mode-service chain, which is rarely seen outside this crew and BlackCat lineage.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1562.009",
        "T1562.001",
        "T1543.003",
        "T1059.001",
        "T1490"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdlines values(Processes.parent_process_name) as parents values(Processes.user) as users from datamodel=Endpoint.Processes where (Processes.process_name IN (\"bcdedit.exe\",\"shutdown.exe\",\"vssadmin.exe\",\"wmic.exe\",\"sc.exe\",\"net.exe\",\"powershell.exe\")) AND (Processes.process=\"*safeboot*\" OR Processes.process=\"*minimal*\" OR Processes.process=\"*network*\" OR Processes.process=\"*shadowcopy delete*\" OR Processes.process=\"*delete shadows*\" OR Processes.process=\"*Set-MpPreference*DisableRealtimeMonitoring*\" OR Processes.process=\"*shutdown*/r*/t*\" OR Processes.process=\"*sc * create *type=*own*start=*\") by host, Processes.process_name, Processes.parent_process_name, Processes.user | `drop_dm_object_name(Processes)` | eventstats dc(process_name) as distinct_actions by host | where distinct_actions>=3 AND match(cmdlines,\"(?i)safeboot|minimal|network\") AND match(cmdlines,\"(?i)/r\\s|/r /t|shutdown\") | sort 0 host firstTime",
      "kql": "let lookback = 6h;\nlet suspiciousProcs = dynamic([\"bcdedit.exe\",\"shutdown.exe\",\"vssadmin.exe\",\"wmic.exe\",\"sc.exe\",\"net.exe\",\"powershell.exe\",\"reg.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(lookback)\n| where FileName in~ (suspiciousProcs)\n| extend Cat = case(\n    FileName =~ \"bcdedit.exe\" and ProcessCommandLine has_any (\"safeboot\",\"minimal\",\"network\"), \"safeboot_set\",\n    FileName =~ \"shutdown.exe\" and ProcessCommandLine has \"/r\", \"forced_reboot\",\n    (FileName =~ \"vssadmin.exe\" or FileName =~ \"wmic.exe\") and ProcessCommandLine has_any (\"delete shadows\",\"shadowcopy delete\"), \"shadow_delete\",\n    FileName =~ \"powershell.exe\" and ProcessCommandLine has_any (\"DisableRealtimeMonitoring\",\"MpPreference\",\"Set-Service\"), \"defender_off\",\n    FileName =~ \"sc.exe\" and ProcessCommandLine has \"create\" and ProcessCommandLine has_any (\"start= auto\",\"type= own\"), \"svc_create\",\n    \"\")\n| where isnotempty(Cat)\n| summarize cats=make_set(Cat), cmds=make_list(ProcessCommandLine,50), parents=make_set(InitiatingProcessFileName), firstSeen=min(Timestamp), lastSeen=max(Timestamp) by DeviceId, DeviceName, bin(Timestamp,30m)\n| where array_length(cats) >= 3 and cats has_any (\"safeboot_set\",\"forced_reboot\")\n| join kind=leftouter (\n    DeviceImageLoadEvents\n    | where Timestamp > ago(lookback)\n    | where FolderPath has_any (\"\\\\probmon.sys\",\"\\\\truesight.sys\",\"\\\\rentdrv2.sys\",\"\\\\gmer.sys\",\"\\\\aswarpot.sys\")\n    | summarize vulnDrivers=make_set(FileName) by DeviceId\n) on DeviceId\n| project firstSeen, lastSeen, DeviceName, cats, vulnDrivers, parents, cmds",
      "rationale": "The combination of a Safe Mode boot flag (bcdedit safeboot minimal/network) plus a forced reboot plus shadow-copy deletion plus Defender tampering plus a service creation inside a single 30-minute window on one host is a near-pathognomonic Embargo / BlackCat-lineage fingerprint. Legitimate admins almost never combine Safe Mode boot manipulation with vssadmin shadow deletion. The optional join agai"
    },
    {
      "title": "[LLM] Embargo Rust loader staging in ProgramData with masquerading service binary",
      "description": "Hunts the file-system staging behavior of Embargo's MDeployer/MS4Killer pair: Rust-compiled PE binaries dropped into C:\\ProgramData (or Public) with names mimicking legitimate utilities (e.g. MDeployer.exe, b.exe, pay.exe, ms4killer.exe), executed by a newly created service or scheduled task, immediately preceding T1486 encryption. This catches the staging step before encryption fires.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "medium",
      "techs": [
        "T1027.009",
        "T1027",
        "T1543.003",
        "T1486",
        "T1489"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Filesystem.file_path) as paths values(Filesystem.process_name) as writers from datamodel=Endpoint.Filesystem where Filesystem.file_path IN (\"*\\\\ProgramData\\\\*\",\"*\\\\Users\\\\Public\\\\*\",\"*\\\\PerfLogs\\\\*\") AND Filesystem.file_name IN (\"MDeployer.exe\",\"ms4killer.exe\",\"MS4Killer.exe\",\"b.exe\",\"pay.exe\",\"praiser.exe\",\"stop.exe\") by host, Filesystem.file_name | `drop_dm_object_name(Filesystem)` | join type=inner host [| tstats summariesonly=true count from datamodel=Endpoint.Processes where Processes.process_name=\"sc.exe\" AND Processes.process=\"*create*\" AND Processes.process IN (\"*ProgramData*\",\"*Users\\\\Public*\") by host Processes.process | `drop_dm_object_name(Processes)`] | sort 0 host firstTime",
      "kql": "let stagingPaths = dynamic([@\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\PerfLogs\\\"]);\nlet badNames = dynamic([\"mdeployer.exe\",\"ms4killer.exe\",\"b.exe\",\"pay.exe\",\"praiser.exe\",\"stop.exe\"]);\nDeviceFileEvents\n| where ActionType in (\"FileCreated\",\"FileRenamed\")\n| where FolderPath has_any (stagingPaths)\n| where FileName in~ (badNames) or (FileName endswith \".exe\" and FileSize between (200000 .. 8000000) and InitiatingProcessFileName in~ (\"powershell.exe\",\"cmd.exe\",\"rundll32.exe\",\"mshta.exe\",\"wmiprvse.exe\"))\n| project DropTime=Timestamp, DeviceId, DeviceName, FileName, FolderPath, SHA256, Dropper=InitiatingProcessFileName, DropperCmd=InitiatingProcessCommandLine\n| join kind=inner (\n    DeviceProcessEvents\n    | where FileName =~ \"sc.exe\" or FileName =~ \"schtasks.exe\"\n    | where ProcessCommandLine has_any (@\"\\ProgramData\\\", @\"\\Users\\Public\\\")\n    | project SvcTime=Timestamp, DeviceId, SvcCmd=ProcessCommandLine\n) on DeviceId\n| where SvcTime between (DropTime .. DropTime + 30m)\n| join kind=leftouter (\n    DeviceProcessEvents\n    | where ProcessCommandLine has_any (\"vssadmin\",\"bcdedit\",\"DisableRealtimeMonitoring\",\"wevtutil cl\")\n    | project ImpactTime=Timestamp, DeviceId, ImpactCmd=ProcessCommandLine\n) on DeviceId\n| where isnull(ImpactTime) or ImpactTime between (DropTime .. DropTime + 1h)\n| project DropTime, DeviceName, FileName, FolderPath, SHA256, Dropper, DropperCmd, SvcCmd, ImpactCmd",
      "rationale": "Embargo's known toolset (MDeployer loader \u2192 MS4Killer EDR-killer \u2192 ransomware payload) consistently lands in writable, non-userprofile directories like ProgramData and is invoked via sc.exe-created services so MS4Killer can run as SYSTEM in Safe Mode. Joining the file drop to a service-creation event referencing the same path within 30 minutes \u2014 and then to a vssadmin/bcdedit/Defender-tamper event"
    }
  ]
}