# intel/.llm_actor_uc_cache/fa/fab83550d470e0a1d267c457b6d989b6026b1868.json
{
  "name": "Lotus Blossom",
  "ucs": [
    {
      "title": "[LLM] Lotus Blossom / Bilbug Sagerunex backdoor C2 to legitimate cloud services from injected svchost.exe",
      "description": "Hunts the signature behaviour of the Sagerunex implant family (Cisco Talos, Mar 2025) used by Lotus Blossom (aka Bilbug / Spring Dragon / Raspberry Typhoon): the implant injects into svchost.exe and beacons to attacker-controlled Dropbox, Twitter (X) and Zimbra mail tenants for C2. Pairing the unusual svchost \u2192 public-cloud-API egress with Sagerunex-style service-registry persistence catches the actor's specific living-off-trusted-infrastructure tradecraft.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "high",
      "techs": [
        "T1543.003",
        "T1102",
        "T1090.001",
        "T1055"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(All_Traffic.dest) as dest values(All_Traffic.dest_port) as dest_port values(All_Traffic.app) as app values(All_Traffic.bytes_out) as bytes_out from datamodel=Network_Traffic.All_Traffic where (All_Traffic.process_name=\"svchost.exe\" OR All_Traffic.app=\"svchost.exe\") AND (All_Traffic.dest_host=\"*dropboxapi.com\" OR All_Traffic.dest_host=\"*content.dropbox*\" OR All_Traffic.dest_host=\"api.twitter.com\" OR All_Traffic.dest_host=\"twitter.com\" OR All_Traffic.dest_host=\"x.com\" OR All_Traffic.dest_host=\"*zimbra*\" OR All_Traffic.url=\"*/oauth2/token*twitter*\" OR All_Traffic.url=\"*/2/files/upload*\") AND NOT (All_Traffic.dest_host IN (\"*.update.microsoft.com\",\"*.windowsupdate.com\",\"*.delivery.mp.microsoft.com\")) by host All_Traffic.user All_Traffic.dest_host All_Traffic.process_name | rename All_Traffic.* as * | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` | join type=outer host [| tstats summariesonly=true count from datamodel=Endpoint.Registry where (Registry.registry_path=\"*\\\\CurrentControlSet\\\\Services\\\\*\" AND Registry.registry_value_name=\"ImagePath\" AND (Registry.registry_value_data=\"*\\\\ProgramData\\\\*\" OR Registry.registry_value_data=\"*\\\\Users\\\\Public\\\\*\" OR Registry.registry_value_data=\"*rundll32*\" OR Registry.registry_value_data=\"*regsvr32*\")) by host Registry.registry_value_data Registry.registry_key_name | rename Registry.* as * | rename count as svc_persist_hits] | where svc_persist_hits >= 1 OR count > 5",
      "kql": "let _ms_update_hosts = dynamic([\"update.microsoft.com\",\"windowsupdate.com\",\"delivery.mp.microsoft.com\",\"download.microsoft.com\",\"sls.update.microsoft.com\"]);\nlet _sagerunex_c2_terms = dynamic([\"dropboxapi.com\",\"content.dropbox\",\"api.twitter.com\",\"twitter.com\",\"x.com/i/api\",\"zimbra\"]);\nlet SagerunexEgress =\n    DeviceNetworkEvents\n    | where Timestamp > ago(14d)\n    | where InitiatingProcessFileName =~ \"svchost.exe\"\n    | where RemoteIPType == \"Public\"\n    | where isnotempty(RemoteUrl)\n    | where RemoteUrl has_any (_sagerunex_c2_terms)\n    | where not(RemoteUrl has_any (_ms_update_hosts))\n    | summarize Beacons = count(),\n                FirstSeen = min(Timestamp),\n                LastSeen  = max(Timestamp),\n                Hosts     = make_set(RemoteUrl, 25),\n                IPs       = make_set(RemoteIP, 25),\n                ParentCmd = any(InitiatingProcessCommandLine)\n                by DeviceId, DeviceName, InitiatingProcessFolderPath\n    | where Beacons >= 3;\nlet SagerunexService =\n    DeviceRegistryEvents\n    | where Timestamp > ago(14d)\n    | where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n    | where RegistryKey has @\"\\CurrentControlSet\\Services\\\"\n    | where RegistryValueName =~ \"ImagePath\" or RegistryValueName =~ \"ServiceDll\"\n    | where RegistryValueData has_any (@\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\AppData\\\", \"rundll32\", \"regsvr32\")\n    | where InitiatingProcessFileName !in~ (\"msiexec.exe\",\"trustedinstaller.exe\",\"services.exe\")\n    | summarize PersistKey = any(RegistryKey),\n                PersistVal = any(RegistryValueData),\n                PersistTime = min(Timestamp)\n                by DeviceId, DeviceName;\nSagerunexEgress\n| join kind=leftouter SagerunexService on DeviceId\n| project FirstSeen, LastSeen, DeviceName, ParentCmd,\n          Beacons, Hosts, IPs,\n          PersistKey, PersistVal, PersistTime\n| order by FirstSeen desc",
      "sentinel_kql": "let _ms_update_hosts = dynamic([\"update.microsoft.com\",\"windowsupdate.com\",\"delivery.mp.microsoft.com\",\"download.microsoft.com\"]);\nlet _sagerunex_c2_terms = dynamic([\"dropboxapi.com\",\"content.dropbox\",\"api.twitter.com\",\"twitter.com\",\"x.com/i/api\",\"zimbra\"]);\nlet SagerunexEgress =\n    CommonSecurityLog\n    | where TimeGenerated > ago(14d)\n    | where SourceProcessName has \"svchost.exe\" or DeviceCustomString1 has \"svchost.exe\" or RequestClientApplication has \"svchost\"\n    | where isnotempty(RequestURL) or isnotempty(DestinationHostName)\n    | extend Target = coalesce(DestinationHostName, RequestURL, DestinationDnsDomain)\n    | where Target has_any (_sagerunex_c2_terms)\n    | where not(Target has_any (_ms_update_hosts))\n    | summarize Beacons = count(),\n                FirstSeen = min(TimeGenerated),\n                LastSeen  = max(TimeGenerated),\n                Targets   = make_set(Target, 25),\n                DstIPs    = make_set(DestinationIP, 25)\n                by SourceHostName, SourceIP\n    | where Beacons >= 3;\nlet SagerunexService =\n    SecurityEvent\n    | where TimeGenerated > ago(14d)\n    | where EventID == 4697 or EventID == 7045\n    | extend SvcImage = coalesce(column_ifexists(\"ServiceFileName\", \"\"), tostring(EventData))\n    | where SvcImage has_any (@\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\AppData\\\", \"rundll32\", \"regsvr32\")\n    | summarize PersistTime = min(TimeGenerated), PersistVal = any(SvcImage) by Computer;\nSagerunexEgress\n| join kind=leftouter SagerunexService on $left.SourceHostName == $right.Computer\n| project FirstSeen, LastSeen, SourceHostName, Beacons, Targets, DstIPs, PersistVal, PersistTime\n| order by FirstSeen desc",
      "sigma_yaml": "title: Lotus Blossom Sagerunex svchost C2 to Dropbox/Twitter/Zimbra\nid: 8e6b3a2c-9f1d-4d9d-9f7a-1c0d4f5b6a91\nstatus: experimental\ndescription: Detects Sagerunex implant beaconing pattern used by Lotus Blossom (Bilbug, Spring Dragon, Raspberry Typhoon) \u2014 svchost.exe egressing to Dropbox API, Twitter/X API, or Zimbra-hosted C2 channels.\nreferences:\n  - https://blog.talosintelligence.com/lotus-blossom-espionage-group/\n  - https://attack.mitre.org/groups/G0030/\nauthor: Detection Engineering\ndate: 2025/03/05\nlogsource:\n  category: network_connection\n  product: windows\ndetection:\n  selection_proc:\n    Image|endswith: '\\svchost.exe'\n  selection_dest:\n    DestinationHostname|contains:\n      - 'dropboxapi.com'\n      - 'content.dropbox'\n      - 'api.twitter.com'\n      - 'zimbra'\n  filter_legit:\n    DestinationHostname|contains:\n      - 'update.microsoft.com'\n      - 'windowsupdate.com'\n      - 'delivery.mp.microsoft.com'\n  condition: selection_proc and selection_dest and not filter_legit\nfields:\n  - Image\n  - DestinationHostname\n  - DestinationIp\n  - User\nfalsepositives:\n  - Genuine corporate Dropbox/Twitter integrations performed via svchost are extremely rare; verify against known third-party drivers/services.\nlevel: high\ntags:\n  - attack.command_and_control\n  - attack.t1102\n  - attack.t1071.001\n  - attack.g0030",
      "rationale": "Sagerunex (a.k.a. Billbug/Bilbug second-stage backdoor evolved from Elise/Evora) is uniquely attributed to Lotus Blossom by Cisco Talos (Mar 2025) and Symantec/Broadcom (2022). Its identifying tradecraft is (a) injection into svchost.exe and (b) C2 over legitimate Dropbox, Twitter and Zimbra services \u2014 abuse Microsoft and other APTs do not commonly mimic. Combining the unusual svchost-to-cloud-API"
    },
    {
      "title": "[LLM] Lotus Blossom on-host recon burst followed by WinRAR staging (Bilbug discovery \u2192 collection chain)",
      "description": "Hunts the multi-step host-recon \u2192 archive-staging chain that Symantec (Bilbug 2022) and Cisco Talos (2025) observed Lotus Blossom operators executing manually after Sagerunex/Hannotog drops a foothold: a tight cluster of net.exe / netstat / ipconfig / tasklist / systeminfo / arp / whoami / nltest from one user, followed within an hour by WinRAR (rar.exe / WinRAR.exe) archiving into ProgramData or Users\\Public.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "medium",
      "techs": [
        "T1016",
        "T1018",
        "T1049",
        "T1057",
        "T1082",
        "T1087.002",
        "T1560.001",
        "T1074.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count values(Processes.process) as cmds values(Processes.process_name) as recon_bins dc(Processes.process_name) as distinct_recon_bins min(_time) as recon_first max(_time) as recon_last from datamodel=Endpoint.Processes where (Processes.process_name IN (\"net.exe\",\"net1.exe\",\"netstat.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"systeminfo.exe\",\"arp.exe\",\"whoami.exe\",\"nltest.exe\",\"route.exe\",\"qwinsta.exe\")) AND NOT Processes.user IN (\"SYSTEM\",\"LOCAL SERVICE\",\"NETWORK SERVICE\",\"*$\") by Processes.user host _time | bin _time span=15m | stats sum(count) as recon_cmd_count values(recon_bins) as recon_bins values(distinct_recon_bins) as distinct_recon_bins min(recon_first) as recon_first max(recon_last) as recon_last by host Processes.user _time | rename Processes.user as user | where distinct_recon_bins >= 5 | join type=inner host user [| tstats summariesonly=true count values(Processes.process) as rar_cmds values(Processes.process_name) as rar_bin min(_time) as archive_time from datamodel=Endpoint.Processes where (Processes.process_name IN (\"rar.exe\",\"WinRAR.exe\",\"7z.exe\") OR Processes.process IN (\"*rar* a *\",\"*7z* a *\")) AND (Processes.process IN (\"*\\\\ProgramData\\\\*\",\"*\\\\Users\\\\Public\\\\*\",\"*\\\\Windows\\\\Temp\\\\*\",\"* -hp*\",\"* -p*\")) by host Processes.user] | where archive_time >= recon_first AND archive_time <= (recon_last + 3600) | table recon_first recon_last archive_time host user distinct_recon_bins recon_bins rar_bin rar_cmds",
      "kql": "let _recon_bins = dynamic([\"net.exe\",\"net1.exe\",\"netstat.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"systeminfo.exe\",\"arp.exe\",\"whoami.exe\",\"nltest.exe\",\"route.exe\",\"qwinsta.exe\"]);\nlet _archive_bins = dynamic([\"rar.exe\",\"winrar.exe\",\"7z.exe\",\"7za.exe\"]);\nlet ReconBurst =\n    DeviceProcessEvents\n    | where Timestamp > ago(14d)\n    | where AccountName !endswith \"$\"\n    | where AccountName !in~ (\"system\",\"local service\",\"network service\")\n    | where FileName in~ (_recon_bins)\n    | summarize ReconCount = count(),\n                DistinctRecon = dcount(FileName),\n                ReconBins = make_set(FileName, 20),\n                SampleCmds = make_set(ProcessCommandLine, 10),\n                ReconStart = min(Timestamp),\n                ReconEnd   = max(Timestamp)\n                by DeviceId, DeviceName, AccountName, bin(Timestamp, 15m)\n    | where DistinctRecon >= 5;          // tight burst of >= 5 distinct discovery binaries in 15m\nlet RarStaging =\n    DeviceProcessEvents\n    | where Timestamp > ago(14d)\n    | where FileName in~ (_archive_bins)\n       or (InitiatingProcessFileName in~ (_archive_bins))\n       or ProcessCommandLine matches regex @\"(?i)\\b(rar|7z)\\b.*\\s+a\\s+\"\n    | where ProcessCommandLine has_any (@\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\Windows\\Temp\\\", \" -hp\", \" -p\", \" -m5\", \" -v\")\n    | project ArchiveTime = Timestamp, DeviceId, DeviceName, AccountName,\n              ArchiveCmd = ProcessCommandLine, ArchiveBin = FileName;\nReconBurst\n| join kind=inner RarStaging on DeviceId, AccountName\n| where ArchiveTime between (ReconStart .. ReconEnd + 1h)\n| extend DelayMin = datetime_diff('minute', ArchiveTime, ReconStart)\n| project ReconStart, ReconEnd, ArchiveTime, DelayMin,\n          DeviceName, AccountName,\n          DistinctRecon, ReconBins, SampleCmds,\n          ArchiveBin, ArchiveCmd\n| order by ReconStart desc",
      "sentinel_kql": "let _recon_bins = dynamic([\"net.exe\",\"net1.exe\",\"netstat.exe\",\"ipconfig.exe\",\"tasklist.exe\",\"systeminfo.exe\",\"arp.exe\",\"whoami.exe\",\"nltest.exe\",\"route.exe\",\"qwinsta.exe\"]);\nlet _archive_bins = dynamic([\"rar.exe\",\"winrar.exe\",\"7z.exe\",\"7za.exe\"]);\nlet ReconBurst =\n    SecurityEvent\n    | where TimeGenerated > ago(14d)\n    | where EventID == 4688\n    | where AccountType == \"User\"\n    | where Account !endswith \"$\"\n    | extend Bin = tolower(tostring(split(NewProcessName, @\"\\\")[-1]))\n    | where Bin in (_recon_bins)\n    | summarize ReconCount = count(),\n                DistinctRecon = dcount(Bin),\n                ReconBins = make_set(Bin, 20),\n                SampleCmds = make_set(CommandLine, 10),\n                ReconStart = min(TimeGenerated),\n                ReconEnd   = max(TimeGenerated)\n                by Computer, Account, bin(TimeGenerated, 15m)\n    | where DistinctRecon >= 5;\nlet RarStaging =\n    SecurityEvent\n    | where TimeGenerated > ago(14d)\n    | where EventID == 4688\n    | extend Bin = tolower(tostring(split(NewProcessName, @\"\\\")[-1]))\n    | where Bin in (_archive_bins) or CommandLine matches regex @\"(?i)\\b(rar|7z)\\b.*\\s+a\\s+\"\n    | where CommandLine has_any (@\"\\ProgramData\\\", @\"\\Users\\Public\\\", @\"\\Windows\\Temp\\\", \" -hp\", \" -p\", \" -m5\", \" -v\")\n    | project ArchiveTime = TimeGenerated, Computer, Account,\n              ArchiveCmd = CommandLine, ArchiveBin = Bin;\nReconBurst\n| join kind=inner RarStaging on Computer, Account\n| where ArchiveTime between (ReconStart .. ReconEnd + 1h)\n| extend DelayMin = datetime_diff('minute', ArchiveTime, ReconStart)\n| project ReconStart, ReconEnd, ArchiveTime, DelayMin,\n          Computer, Account, DistinctRecon, ReconBins, SampleCmds,\n          ArchiveBin, ArchiveCmd\n| order by ReconStart desc",
      "sigma_yaml": "",
      "rationale": "Symantec's 2022 Bilbug report and Talos's 2025 Lotus Blossom report both call out the same hands-on-keyboard recon-then-RAR pattern as a Lotus Blossom signature: operators run a tight burst of native discovery commands (net/netstat/ipconfig/tasklist/systeminfo/arp/whoami/nltest) and then archive findings with WinRAR (often password-protected with -hp) into ProgramData or Users\\Public for later exf"
    }
  ]
}