# intel/.llm_actor_uc_cache/30/30c81995cb3a9749215bb23fb99f662d4e9488f8.json
{
  "name": "Lazarus Group",
  "ucs": [
    {
      "title": "[LLM] Lazarus 'Operation DreamJob' / 3CX-style trojanized signed binary side-loading a malicious DLL from a user-writable path",
      "description": "Hunts the canonical Lazarus initial-access chain (Diamond Sleet / Labyrinth Chollima): a job-themed ISO/ZIP drops a renamed copy of a legitimate signed application (SumatraPDF, CiscoCollabHost, FoxitReader, X_TRADER, 3CXDesktopApp, etc.) into Public/Temp/Downloads/ProgramData; that binary loads an adversary-supplied side-load DLL (LightlessCan, MagicRAT, BLINDINGCAN, AppleJeus loaders) from the same directory. The pairing of a Lazarus-trojanized application name + adjacent side-load DLL + user-writable path is what separates this from generic side-loading.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574",
        "T1036.005",
        "T1027",
        "T1027.009",
        "T1059.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.process_name IN (\"SumatraPDF.exe\",\"CiscoCollabHost.exe\",\"colorui.exe\",\"PCKeeper.exe\",\"FoxitReader.exe\",\"X_TRADER.exe\",\"3CXDesktopApp.exe\",\"ZoomLauncher.exe\",\"WSearch.exe\",\"GoogleUpdate.exe\") (Processes.process_path=\"*\\\\Users\\\\Public\\\\*\" OR Processes.process_path=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Processes.process_path=\"*\\\\Downloads\\\\*\" OR Processes.process_path=\"*\\\\ProgramData\\\\*\") NOT (Processes.user IN (\"*$\",\"SYSTEM\",\"NETWORK SERVICE\",\"LOCAL SERVICE\")) by Processes.dest Processes.user Processes.process_name Processes.process Processes.process_path Processes.parent_process_name Processes.parent_process Processes.process_hash | `drop_dm_object_name(Processes)` | rename process_path as launcher_path, process_hash as launcher_hash | join type=inner dest launcher_path [ search (sourcetype=\"XmlWinEventLog:Microsoft-Windows-Sysmon/Operational\" EventCode=7) (ImageLoaded=\"*\\\\msvcr100.dll\" OR ImageLoaded=\"*\\\\libcef.dll\" OR ImageLoaded=\"*\\\\colorui.dll\" OR ImageLoaded=\"*\\\\wsock32.dll\" OR ImageLoaded=\"*\\\\VERSION.dll\" OR ImageLoaded=\"*\\\\msvcp140.dll\" OR ImageLoaded=\"*\\\\HID.DLL\" OR ImageLoaded=\"*\\\\ssleay32.dll\" OR ImageLoaded=\"*\\\\libeay32.dll\" OR ImageLoaded=\"*\\\\wlbsctrl.dll\") | rex field=Image \"(?<launcher_path>.*\\\\\\\\)\" | rename Computer as dest | fields dest launcher_path ImageLoaded ImageLoadedHash Signed ] | where Signed!=\"true\" OR isnull(Signed) | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)` | table firstTime lastTime dest user process_name launcher_path launcher_hash parent_process_name ImageLoaded ImageLoadedHash",
      "kql": "let _LazarusTrojanizedBins = dynamic([\"SumatraPDF.exe\",\"CiscoCollabHost.exe\",\"colorui.exe\",\"PCKeeper.exe\",\"FoxitReader.exe\",\"X_TRADER.exe\",\"3CXDesktopApp.exe\",\"ZoomLauncher.exe\",\"GoogleUpdate.exe\",\"WSearch.exe\"]);\nlet _LazarusSideLoadedDlls = dynamic([\"msvcr100.dll\",\"libcef.dll\",\"colorui.dll\",\"HID.DLL\",\"ssleay32.dll\",\"libeay32.dll\",\"wlbsctrl.dll\",\"VERSION.dll\",\"wsock32.dll\",\"winmm.dll\",\"msvcp140.dll\",\"secur32.dll\",\"DSROLE.dll\"]);\nlet _userPaths = dynamic([@\"\\Users\\Public\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Downloads\\\", @\"\\ProgramData\\\", @\"\\AppData\\Roaming\\\"]);\nlet Loaders = DeviceProcessEvents\n    | where Timestamp > ago(14d)\n    | where InitiatingProcessAccountName !endswith \"$\"\n    | where FileName in~ (_LazarusTrojanizedBins)\n    | where FolderPath has_any (_userPaths)\n    | project Timestamp, DeviceId, DeviceName, AccountName,\n              LoaderProcessId = ProcessId,\n              LoaderName = FileName, LoaderPath = FolderPath,\n              LoaderCmd = ProcessCommandLine, LoaderSHA256 = SHA256,\n              LoaderCompany = ProcessVersionInfoCompanyName,\n              ParentName = InitiatingProcessFileName,\n              ParentCmd = InitiatingProcessCommandLine;\nDeviceImageLoadEvents\n| where Timestamp > ago(14d)\n| where FileName in~ (_LazarusSideLoadedDlls)\n| where FolderPath has_any (_userPaths)\n| join kind=inner Loaders on DeviceId, $left.InitiatingProcessId == $right.LoaderProcessId\n| where FolderPath =~ LoaderPath                           // DLL adjacent to loader\n| where SHA256 != LoaderSHA256                              // not the loader itself\n| where isempty(InitiatingProcessVersionInfoCompanyName)    // unsigned / blank company\n     or InitiatingProcessVersionInfoCompanyName !~ LoaderCompany   // company mismatch\n| project Timestamp, DeviceName, AccountName, ParentName, ParentCmd,\n          LoaderName, LoaderPath, LoaderCmd, LoaderSHA256, LoaderCompany,\n          SideLoadedDll = FileName, DllPath = FolderPath, DllSHA256 = SHA256,\n          DllCompany = InitiatingProcessVersionInfoCompanyName\n| order by Timestamp desc",
      "sentinel_kql": "let _LazarusTrojanizedBins = dynamic([\"SumatraPDF.exe\",\"CiscoCollabHost.exe\",\"colorui.exe\",\"PCKeeper.exe\",\"FoxitReader.exe\",\"X_TRADER.exe\",\"3CXDesktopApp.exe\",\"ZoomLauncher.exe\",\"GoogleUpdate.exe\",\"WSearch.exe\"]);\nlet _LazarusSideLoadedDlls = dynamic([\"msvcr100.dll\",\"libcef.dll\",\"colorui.dll\",\"HID.DLL\",\"ssleay32.dll\",\"libeay32.dll\",\"wlbsctrl.dll\",\"VERSION.dll\",\"wsock32.dll\",\"winmm.dll\",\"msvcp140.dll\",\"secur32.dll\"]);\nlet _userPaths = dynamic([@\"\\Users\\Public\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Downloads\\\", @\"\\ProgramData\\\", @\"\\AppData\\Roaming\\\"]);\nlet Loaders = SecurityEvent\n    | where TimeGenerated > ago(14d)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ProcFile = tostring(split(NewProcessName, @\"\\\")[-1])\n    | where ProcFile in~ (_LazarusTrojanizedBins)\n    | where NewProcessName has_any (_userPaths)\n    | project LoaderTime = TimeGenerated, Computer, Account,\n              LoaderProcId = NewProcessId, LoaderImage = NewProcessName,\n              LoaderCmd = CommandLine;\nlet Sysmon7 = WindowsEvent\n    | where TimeGenerated > ago(14d)\n    | where Provider == \"Microsoft-Windows-Sysmon\" and EventID == 7\n    | extend Image = tostring(EventData.Image),\n             ImageLoaded = tostring(EventData.ImageLoaded),\n             Signed = tostring(EventData.Signed),\n             Signature = tostring(EventData.Signature),\n             ProcessId = tostring(EventData.ProcessId),\n             Hashes = tostring(EventData.Hashes)\n    | extend DllFile = tostring(split(ImageLoaded, @\"\\\")[-1])\n    | where DllFile in~ (_LazarusSideLoadedDlls)\n    | where ImageLoaded has_any (_userPaths);\nSysmon7\n| join kind=inner Loaders on $left.Computer == $right.Computer, $left.Image == $right.LoaderImage\n| where Signed =~ \"false\" or Signature has_any (\"Invalid\",\"Unavailable\",\"-\")\n| project TimeGenerated, Computer, Account, LoaderImage, LoaderCmd,\n          ImageLoaded, Signed, Signature, Hashes\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "This is the textbook Lazarus tradecraft documented across Operation DreamJob (Mandiant), Operation Interception (ESET), the 3CX supply-chain attack (Microsoft DEV-0139/Diamond Sleet), the X_TRADER compromise, and AppleJeus/MagicRAT/LightlessCan campaigns. Lazarus consistently pairs a renamed-but-still-validly-signed legitimate application with a malicious side-load DLL named after a real Windows/C"
    },
    {
      "title": "[LLM] Lazarus rundll32/regsvr32 named-export persistence pointed at a DLL in %ProgramData% / %Public% via scheduled task or service",
      "description": "Hunts the post-install persistence stage that Diamond Sleet / Labyrinth Chollima uses across BLINDINGCAN, MagicRAT, COPPERHEDGE and the 3CX second-stage implant: after the side-loaded loader writes a payload DLL into ProgramData/Public/AppData, persistence is established as a scheduled task or auto-start service that invokes `rundll32.exe <path>,<NamedExport>` (commonly export ordinals or names like `Init`, `RtlInit`, `_StartW`). The rundll32 grandparent-chain (services.exe / svchost.exe / taskeng.exe) plus a non-system DLL path is the high-fidelity pivot.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1053.005",
        "T1218.011",
        "T1036.005",
        "T1027",
        "T1071.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process) as parent_cmdline values(Processes.process_hash) as process_hash from datamodel=Endpoint.Processes where Processes.process_name IN (\"rundll32.exe\",\"regsvr32.exe\") (Processes.process=\"*\\\\Users\\\\Public\\\\*.dll*\" OR Processes.process=\"*\\\\ProgramData\\\\*.dll*\" OR Processes.process=\"*\\\\AppData\\\\Roaming\\\\*.dll*\" OR Processes.process=\"*\\\\AppData\\\\Local\\\\Temp\\\\*.dll*\") Processes.parent_process_name IN (\"svchost.exe\",\"taskeng.exe\",\"schtasks.exe\",\"services.exe\",\"wmiprvse.exe\") NOT Processes.user IN (\"*$\") by Processes.dest Processes.user Processes.process_name Processes.parent_process_name | `drop_dm_object_name(Processes)` | rex field=cmdline \"(?i)\\\"?(?<dll_path>[a-z]:\\\\\\\\[^\\\",]+\\.dll)\\\"?[\\s,]+(?<export>[A-Za-z_#@][\\w#@]{0,40})\" | where isnotnull(export) AND NOT match(dll_path, \"(?i)\\\\\\\\Windows\\\\\\\\|\\\\\\\\Program Files\") | stats min(firstTime) as firstTime max(lastTime) as lastTime values(cmdline) as cmdline values(dll_path) as dll_path values(export) as export values(parent_cmdline) as parent_cmdline values(process_hash) as process_hash by dest user process_name parent_process_name | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "DeviceProcessEvents\n| where Timestamp > ago(30d)\n| where AccountName !endswith \"$\"\n| where FileName in~ (\"rundll32.exe\",\"regsvr32.exe\")\n| where InitiatingProcessFileName in~ (\"svchost.exe\",\"taskeng.exe\",\"schtasks.exe\",\"services.exe\",\"wmiprvse.exe\")\n| where ProcessCommandLine matches regex @'(?i)([a-z]:\\\\(?:Users\\\\Public|ProgramData|Users\\\\[^\\\\]+\\\\AppData\\\\(?:Roaming|Local\\\\Temp))\\\\[^\\s,\"]+\\.dll)[\\s\",]+([A-Za-z_#@][\\w#@]{0,40})'\n| extend DllPath  = extract(@'(?i)([a-z]:\\\\(?:Users\\\\Public|ProgramData|Users\\\\[^\\\\]+\\\\AppData\\\\(?:Roaming|Local\\\\Temp))\\\\[^\\s,\"]+\\.dll)', 1, ProcessCommandLine)\n| extend Export   = extract(@'(?i)\\.dll[\\s\",]+([A-Za-z_#@][\\w#@]{0,40})', 1, ProcessCommandLine)\n| where isnotempty(DllPath) and isnotempty(Export)\n| join kind=leftouter (\n    DeviceFileEvents\n    | where Timestamp > ago(30d)\n    | where ActionType in (\"FileCreated\",\"FileModified\")\n    | where FileName endswith \".dll\"\n    | project DllSHA256 = SHA256, DllFolderPath = FolderPath, DllFileName = FileName,\n              DllCreatorAccount = InitiatingProcessAccountName,\n              DllCreatorName = InitiatingProcessFileName,\n              DllCreateTime = Timestamp, DeviceId\n  ) on DeviceId\n| where DllPath has DllFileName and DllPath has DllFolderPath\n| project Timestamp, DeviceName, AccountName,\n          GrandParent = InitiatingProcessParentFileName,\n          Parent = InitiatingProcessFileName,\n          ParentCmd = InitiatingProcessCommandLine,\n          RundllCmd = ProcessCommandLine,\n          DllPath, Export, DllSHA256, DllCreatorName, DllCreatorAccount, DllCreateTime\n| order by Timestamp desc",
      "sentinel_kql": "SecurityEvent\n| where TimeGenerated > ago(30d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| extend Image  = tolower(tostring(split(NewProcessName, @\"\\\")[-1])),\n         Parent = tolower(tostring(split(ParentProcessName, @\"\\\")[-1]))\n| where Image in (\"rundll32.exe\",\"regsvr32.exe\")\n| where Parent in (\"svchost.exe\",\"taskeng.exe\",\"schtasks.exe\",\"services.exe\",\"wmiprvse.exe\")\n| where CommandLine matches regex @'(?i)([a-z]:\\\\(?:Users\\\\Public|ProgramData|Users\\\\[^\\\\]+\\\\AppData\\\\(?:Roaming|Local\\\\Temp))\\\\[^\\s,\"]+\\.dll)[\\s\",]+([A-Za-z_#@][\\w#@]{0,40})'\n| extend DllPath = extract(@'(?i)([a-z]:\\\\(?:Users\\\\Public|ProgramData|Users\\\\[^\\\\]+\\\\AppData\\\\(?:Roaming|Local\\\\Temp))\\\\[^\\s,\"]+\\.dll)', 1, CommandLine),\n         Export  = extract(@'(?i)\\.dll[\\s\",]+([A-Za-z_#@][\\w#@]{0,40})', 1, CommandLine)\n| where isnotempty(DllPath) and isnotempty(Export)\n| project TimeGenerated, Computer, Account, ParentProcessName, NewProcessName,\n          CommandLine, DllPath, Export\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Lazarus rundll32/regsvr32 Named-Export Persistence from User-Writable Path\nid: 9b2e2a7e-3f4a-4f6c-9c4d-d9c7c0a6f1a4\nstatus: experimental\ndescription: Detects rundll32.exe or regsvr32.exe spawned by a service/scheduler parent and pointed at a DLL with an explicit named export located in ProgramData/Public/AppData. Matches the Diamond Sleet / Labyrinth Chollima (Lazarus) post-install persistence pattern observed across BLINDINGCAN, MagicRAT, COPPERHEDGE and 3CX second-stage implants.\nreferences:\n  - https://www.microsoft.com/en-us/security/blog/2023/04/20/3cx-supply-chain-attack/\n  - https://attack.mitre.org/groups/G0032/\nauthor: Detection Engineering\ndate: 2026-05-05\ntags:\n  - attack.persistence\n  - attack.defense_evasion\n  - attack.t1053.005\n  - attack.t1218.011\n  - attack.g0032\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_image:\n    Image|endswith:\n      - '\\rundll32.exe'\n      - '\\regsvr32.exe'\n  selection_parent:\n    ParentImage|endswith:\n      - '\\svchost.exe'\n      - '\\taskeng.exe'\n      - '\\schtasks.exe'\n      - '\\services.exe'\n      - '\\wmiprvse.exe'\n  selection_path:\n    CommandLine|re|i: '([a-z]:\\\\(?:Users\\\\Public|ProgramData|Users\\\\[^\\\\]+\\\\AppData\\\\(?:Roaming|Local\\\\Temp))\\\\[^\\s,\"]+\\.dll)[\\s\",]+([A-Za-z_#@][\\w#@]{0,40})'\n  filter_user:\n    User|endswith: '$'\n  condition: selection_image and selection_parent and selection_path and not filter_user\nfields:\n  - Image\n  - ParentImage\n  - CommandLine\n  - User\n  - Computer\nfalsepositives:\n  - Legitimate enterprise software persisting via scheduled task pointing at ProgramData (rare; profile per-org)\n  - Microsoft Office click-to-run repair tasks invoking rundll32 against Office paths\nlevel: high\n",
      "rationale": "Lazarus's mid-stage tooling (BLINDINGCAN, MagicRAT, COPPERHEDGE, 3CX second-stage `iconic.dll`) is consistently invoked via rundll32/regsvr32 with an explicit export name against a payload DLL written to ProgramData or Public \u2014 Mandiant, Microsoft (Diamond Sleet) and Cisco Talos all document this exact pattern. The combination of (a) service/scheduler grandparent, (b) DLL path in user-writable spa"
    }
  ]
}