# intel/.llm_actor_uc_cache/c3/c3ac6f1da27e9fb8ba4bf2ebc8a8634a80f1ff32.json
{
  "name": "APT41",
  "ucs": [
    {
      "title": "[LLM] APT41 (Brass Typhoon / Wicked Panda) web-app exploitation: webserver process spawning recon LOLBins",
      "description": "Detects post-exploitation reconnaissance commands spawned directly by internet-facing webserver worker processes (IIS w3wp.exe, Tomcat java.exe, ManageEngine bundled JRE). APT41 is the most prolific abuser of n-day flaws in internet-exposed apps \u2014 Citrix NetScaler ADC (CVE-2019-19781), Zoho ManageEngine ADSelfService Plus / ServiceDesk (CVE-2021-40539, CVE-2022-47966), Atlassian Confluence (CVE-2022-26134), and Log4Shell \u2014 and characteristically drops a JSP/ASPX webshell that immediately runs whoami / net / nltest / systeminfo to enumerate the foothold before deploying DUSTPAN / DUSTTRAP / KEY",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "exploit",
      "conf": "high",
      "techs": [
        "T1190",
        "T1059.001",
        "T1059.003",
        "T1087.002",
        "T1082",
        "T1018"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as command_lines values(Processes.process_name) as child_procs dc(Processes.process_name) as distinct_children from datamodel=Endpoint.Processes where (Processes.parent_process_name IN (\"w3wp.exe\",\"tomcat.exe\",\"tomcat9.exe\",\"tomcat10.exe\",\"java.exe\",\"javaw.exe\",\"httpd.exe\",\"nginx.exe\")) AND (Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"net.exe\",\"net1.exe\",\"whoami.exe\",\"ipconfig.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"nltest.exe\",\"certutil.exe\",\"quser.exe\",\"hostname.exe\",\"reg.exe\")) by Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process\n| `drop_dm_object_name(Processes)`\n| where NOT match(user, \"\\$$\")\n| convert ctime(firstTime) ctime(lastTime)\n| sort - lastTime",
      "kql": "// APT41 webshell post-ex recon \u2014 webserver worker spawning native discovery LOLBins\nlet WebServerParents = dynamic([\"w3wp.exe\",\"tomcat.exe\",\"tomcat9.exe\",\"tomcat10.exe\",\"java.exe\",\"javaw.exe\",\"httpd.exe\",\"nginx.exe\"]);\nlet ReconBins = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"net.exe\",\"net1.exe\",\"whoami.exe\",\"ipconfig.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"nltest.exe\",\"certutil.exe\",\"quser.exe\",\"hostname.exe\",\"reg.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (WebServerParents)\n| where FileName in~ (ReconBins)\n| where AccountName !endswith \"$\"            // exclude SYSTEM / machine accounts\n| project Timestamp, DeviceName, AccountName,\n          Parent = InitiatingProcessFileName,\n          ParentPath = InitiatingProcessFolderPath,\n          ParentCmd = InitiatingProcessCommandLine,\n          Child = FileName,\n          ChildCmd = ProcessCommandLine,\n          IntegrityLevel = ProcessIntegrityLevel,\n          SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// APT41 webshell post-ex recon \u2014 Sentinel via SecurityEvent 4688\nlet WebServerParents = dynamic([@\"\\w3wp.exe\",@\"\\tomcat.exe\",@\"\\tomcat9.exe\",@\"\\tomcat10.exe\",@\"\\java.exe\",@\"\\javaw.exe\",@\"\\httpd.exe\",@\"\\nginx.exe\"]);\nlet ReconBins = dynamic([@\"\\cmd.exe\",@\"\\powershell.exe\",@\"\\pwsh.exe\",@\"\\net.exe\",@\"\\net1.exe\",@\"\\whoami.exe\",@\"\\ipconfig.exe\",@\"\\systeminfo.exe\",@\"\\tasklist.exe\",@\"\\nltest.exe\",@\"\\certutil.exe\",@\"\\quser.exe\",@\"\\hostname.exe\",@\"\\reg.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where ParentProcessName has_any (WebServerParents)\n| where NewProcessName has_any (ReconBins)\n| where Account !endswith \"$\"\n| project TimeGenerated, Computer, Account,\n          ParentProcessName, NewProcessName, CommandLine,\n          TokenElevationType, SubjectUserName\n| order by TimeGenerated desc",
      "sigma_yaml": "title: APT41 Web Application Exploitation Followed by Webshell Recon\nid: a4d3b1e8-7c2f-4a5b-9e8d-6f1c2b3a4e5f\nstatus: experimental\ndescription: |\n  Detects native discovery LOLBins spawned directly by internet-facing webserver\n  worker processes (IIS w3wp, Tomcat java, Apache, nginx). Consistent with APT41 /\n  Brass Typhoon / Wicked Panda webshell deployment after exploitation of\n  Confluence (CVE-2022-26134), Log4Shell, ManageEngine (CVE-2021-40539,\n  CVE-2022-47966), and Citrix NetScaler ADC (CVE-2019-19781).\nreferences:\n  - https://attack.mitre.org/groups/G0096/\n  - https://www.mandiant.com/resources/blog/apt41-us-state-governments\nauthor: Detection Engineering\ndate: 2026/05/06\ntags:\n  - attack.initial_access\n  - attack.t1190\n  - attack.execution\n  - attack.t1059.001\n  - attack.t1059.003\n  - attack.discovery\n  - attack.t1087.002\n  - attack.g0096\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent:\n    ParentImage|endswith:\n      - '\\w3wp.exe'\n      - '\\tomcat.exe'\n      - '\\tomcat9.exe'\n      - '\\tomcat10.exe'\n      - '\\java.exe'\n      - '\\javaw.exe'\n      - '\\httpd.exe'\n      - '\\nginx.exe'\n  selection_child:\n    Image|endswith:\n      - '\\cmd.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\net.exe'\n      - '\\net1.exe'\n      - '\\whoami.exe'\n      - '\\ipconfig.exe'\n      - '\\systeminfo.exe'\n      - '\\tasklist.exe'\n      - '\\nltest.exe'\n      - '\\certutil.exe'\n      - '\\quser.exe'\n      - '\\hostname.exe'\n      - '\\reg.exe'\n  filter_machine:\n    User|endswith: '$'\n  condition: selection_parent and selection_child and not filter_machine\nfalsepositives:\n  - Legitimate site administration scripts triggered via web management consoles\n  - APM / monitoring agents probing the host\nlevel: high",
      "rationale": "Webserver workers (w3wp, tomcat java) almost never have a legitimate reason to spawn cmd/powershell/whoami/nltest. APT41 has used this exact spawn pattern across multiple campaigns documented by Mandiant (US state governments 2022, ProxyLogon piggyback, ManageEngine ADSelfService Plus exploitation), and CrowdStrike's Wicked Panda profile. The parent restriction strips ~99% of legit recon noise; pa"
    },
    {
      "title": "[LLM] APT41 certutil staging chained to schtasks persistence (DUSTPAN / StealthVector loader pattern)",
      "description": "Detects the canonical APT41 / Wicked Panda staging chain: certutil.exe used as a download/decoder tool followed within ~1 hour by schtasks.exe creating a persistence task on the same host. APT41 has consistently used certutil for retrieving/base64-decoding the StealthVector (DUSTPAN) loader and KEYPLUG / Cobalt Strike payloads, then registered them as scheduled tasks (often masquerading under benign-looking service names) for boot persistence \u2014 see Mandiant M-Trends 2024 and Google TAG / Mandiant DUSTTRAP reporting.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1105",
        "T1027.013",
        "T1053.005",
        "T1036.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t min(_time) as firstTime max(_time) as lastTime values(Processes.process) as commands values(Processes.process_name) as proc_names dc(Processes.process_name) as nproc from datamodel=Endpoint.Processes where (Processes.process_name=\"certutil.exe\" AND (Processes.process=\"*-urlcache*\" OR Processes.process=\"*-decode*\" OR Processes.process=\"*-decodehex*\" OR Processes.process=\"*-split*\" OR Processes.process=\"*/urlcache*\" OR Processes.process=\"*/decode*\")) OR (Processes.process_name=\"schtasks.exe\" AND (Processes.process=\"*/create*\" OR Processes.process=\"*-create*\")) by Processes.dest Processes.user\n| `drop_dm_object_name(Processes)`\n| where nproc>=2 AND (lastTime-firstTime)<=3600 AND NOT match(user, \"\\$$\")\n| eval delay_seconds=(lastTime-firstTime)\n| convert ctime(firstTime) ctime(lastTime)\n| sort - lastTime",
      "kql": "// APT41 staging chain: certutil download/decode followed by schtasks persistence within 1h on the same host\nlet LookbackHours = 24h;\nlet CorrelationWindow = 1h;\nlet CertutilStaging =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackHours)\n    | where FileName =~ \"certutil.exe\"\n    | where ProcessCommandLine has_any (\"-urlcache\",\"/urlcache\",\"-decode\",\"/decode\",\"-decodehex\",\"-split\",\"-f http\")\n    | where AccountName !endswith \"$\"\n    | project StagingTime = Timestamp, DeviceId, DeviceName, AccountName,\n              StagingCmd = ProcessCommandLine,\n              StagingParent = InitiatingProcessFileName;\nlet TaskCreate =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackHours)\n    | where FileName =~ \"schtasks.exe\"\n    | where ProcessCommandLine has_any (\"/create\",\"-create\")\n    | project TaskTime = Timestamp, DeviceId, DeviceName, AccountName,\n              TaskCmd = ProcessCommandLine,\n              TaskParent = InitiatingProcessFileName;\nCertutilStaging\n| join kind=inner TaskCreate on DeviceId\n| where TaskTime between (StagingTime .. StagingTime + CorrelationWindow)\n| project DeviceName, AccountName,\n          StagingTime, TaskTime,\n          DelayMin = datetime_diff('minute', TaskTime, StagingTime),\n          StagingParent, StagingCmd,\n          TaskParent, TaskCmd\n| order by StagingTime desc",
      "sentinel_kql": "// APT41 certutil staging \u2192 schtasks persistence (Sentinel SecurityEvent 4688)\nlet LookbackHours = 24h;\nlet CorrelationWindow = 1h;\nlet CertutilStaging =\n    SecurityEvent\n    | where TimeGenerated > ago(LookbackHours)\n    | where EventID == 4688\n    | where NewProcessName endswith @\"\\certutil.exe\"\n    | where CommandLine has_any (\"-urlcache\",\"/urlcache\",\"-decode\",\"/decode\",\"-decodehex\",\"-split\",\"-f http\")\n    | where Account !endswith \"$\"\n    | project StagingTime = TimeGenerated, Computer, Account,\n              StagingCmd = CommandLine,\n              StagingParent = ParentProcessName;\nlet TaskCreate =\n    SecurityEvent\n    | where TimeGenerated > ago(LookbackHours)\n    | where EventID == 4688\n    | where NewProcessName endswith @\"\\schtasks.exe\"\n    | where CommandLine has_any (\"/create\",\"-create\")\n    | project TaskTime = TimeGenerated, Computer, Account,\n              TaskCmd = CommandLine,\n              TaskParent = ParentProcessName;\nCertutilStaging\n| join kind=inner TaskCreate on Computer\n| where TaskTime between (StagingTime .. StagingTime + CorrelationWindow)\n| project Computer, Account,\n          StagingTime, TaskTime,\n          DelayMin = datetime_diff('minute', TaskTime, StagingTime),\n          StagingParent, StagingCmd,\n          TaskParent, TaskCmd\n| order by StagingTime desc",
      "sigma_yaml": "",
      "rationale": "certutil.exe with -urlcache/-decode is rarely used by legit IT (BITS, MSI, or signed installers handle download/decode), and APT41 has been documented across the FBI/CISA Aug-2022 advisory, Mandiant's APT41 US-state campaign report, and Mandiant's M-Trends 2024 DUSTTRAP write-up using exactly this two-step pattern: certutil pulls/decodes the StealthVector (DUSTPAN) loader, then schtasks pins it to"
    }
  ]
}