# intel/.llm_actor_uc_cache/dd/dd4649d5fd55319e95f646224d1b5a644f5e9e70.json
{
  "name": "Saint Bear",
  "ucs": [
    {
      "title": "[LLM] Saint Bear (UAC-0056) Saint Bot / OutSteel delivery \u2014 Office or browser spawns script interpreter that drops PE to %AppData%/%Temp%",
      "description": "Detects the canonical Saint Bear (Lorec53 / TA471 / Storm-0587) initial-access-to-install chain: a Ukrainian-government-themed phishing attachment or web-staged ZIP/SFX is opened, Office/Explorer/browser spawns wscript/cscript/mshta/powershell, which in turn writes a Saint Bot loader or OutSteel stealer binary into %AppData%\\Roaming, %AppData%\\Local\\Temp, or %Users%\\Public within seconds. This is the exact tradecraft Unit 42, CERT-UA, and Cybereason document for UAC-0056 since 2021.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.007",
        "T1059.001",
        "T1059.003",
        "T1027.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstSeen max(_time) as lastSeen from datamodel=Endpoint.Filesystem where Filesystem.action=created AND (Filesystem.file_path=\"*\\\\AppData\\\\Roaming\\\\*\" OR Filesystem.file_path=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Filesystem.file_path=\"*\\\\Users\\\\Public\\\\*\") AND (Filesystem.file_name=\"*.exe\" OR Filesystem.file_name=\"*.dll\" OR Filesystem.file_name=\"*.scr\" OR Filesystem.file_name=\"*.cpl\") AND Filesystem.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\") by host Filesystem.dest Filesystem.process_name Filesystem.process_guid Filesystem.file_name Filesystem.file_path Filesystem.user | `drop_dm_object_name(Filesystem)` | join type=inner dest process_guid [ | tstats summariesonly=true count from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"explorer.exe\",\"7zfm.exe\",\"winrar.exe\") AND Processes.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\") by Processes.dest Processes.process_guid Processes.process_name Processes.parent_process_name Processes.process | `drop_dm_object_name(Processes)` ] | where firstSeen >= relative_time(now(),\"-7d\") | table firstSeen dest user parent_process_name process_name process file_name file_path",
      "kql": "let win = 180s;\nlet ScriptFromUser = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where AccountName !endswith \"$\"\n| where InitiatingProcessFileName in~ (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"explorer.exe\",\"7zfm.exe\",\"winrar.exe\")\n| where FileName in~ (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\")\n| project ScriptTime = Timestamp, DeviceId, DeviceName, AccountName,\n          ScriptParent = InitiatingProcessFileName,\n          ScriptFile  = FileName,\n          ScriptCmd   = ProcessCommandLine;\nDeviceFileEvents\n| where Timestamp > ago(7d)\n| where ActionType == \"FileCreated\"\n| where FolderPath has_any (@\"\\AppData\\Roaming\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Users\\Public\\\")  // Saint Bot / OutSteel drop sites\n| where FileName endswith \".exe\" or FileName endswith \".dll\"\n     or FileName endswith \".scr\" or FileName endswith \".cpl\"\n| where InitiatingProcessFileName in~ (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\")\n| join kind=inner ScriptFromUser on DeviceId\n| where Timestamp between (ScriptTime .. ScriptTime + win)\n| project Timestamp, DeviceName, AccountName,\n          DroppedFile = FileName, FolderPath, SHA256, FileOriginUrl,\n          ScriptParent, ScriptFile, ScriptCmd,\n          DelaySec = datetime_diff('second', Timestamp, ScriptTime)\n| order by Timestamp desc",
      "sentinel_kql": "let win = 180s;\nlet ScriptFromUser = ImProcessCreate\n| where TimeGenerated > ago(7d)\n| where ActorUsername !endswith \"$\"\n| where ParentProcessName has_any (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"explorer.exe\",\"7zfm.exe\",\"winrar.exe\")\n| where TargetProcessName has_any (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\")\n| project ScriptTime = TimeGenerated, DvcHostname, ActorUsername,\n          ScriptParent = ParentProcessName,\n          ScriptFile   = TargetProcessName,\n          ScriptCmd    = TargetProcessCommandLine;\nImFileEvent\n| where TimeGenerated > ago(7d)\n| where EventType == \"FileCreated\"\n| where TargetFilePath has_any (@\"\\AppData\\Roaming\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Users\\Public\\\")\n| where TargetFileName endswith \".exe\" or TargetFileName endswith \".dll\"\n     or TargetFileName endswith \".scr\" or TargetFileName endswith \".cpl\"\n| where ActingProcessName has_any (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\")\n| join kind=inner ScriptFromUser on DvcHostname\n| where TimeGenerated between (ScriptTime .. ScriptTime + win)\n| project TimeGenerated, DvcHostname, ActorUsername,\n          DroppedFile = TargetFileName, TargetFilePath, TargetFileSHA256,\n          ScriptParent, ScriptFile, ScriptCmd,\n          DelaySec = datetime_diff('second', TimeGenerated, ScriptTime)\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "Saint Bear's defining initial-access pattern (Unit 42 'OutSteel / Saint Bot' report, CERT-UA UAC-0056 alerts) is a Ukrainian-government-spoofed phish carrying a script (.js/.vbs/.lnk/.hta) inside an SFX/ZIP that, when executed, uses wscript/cscript/mshta/powershell to fetch and stage the Saint Bot loader or OutSteel stealer in %AppData% or %Temp%. Joining Office/Explorer/browser \u2192 script-interpret"
    },
    {
      "title": "[LLM] Saint Bot Run-key persistence \u2014 script interpreter writes %AppData% binary path to HKCU\\...\\Run",
      "description": "Hunts Saint Bear's Saint Bot persistence step: after the dropper unpacks, the same script interpreter (wscript/cscript/cmd/powershell/rundll32) immediately registers an autorun under HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run (or RunOnce) whose value points back into %AppData%\\Roaming, %AppData%\\Local\\Temp, or %Users%\\Public. This combination is rare among legitimate installers (which use MSI/Setup) and is documented across Saint Bot samples since 2021.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1547.001",
        "T1112",
        "T1059.001",
        "T1059.007"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstSeen max(_time) as lastSeen from datamodel=Endpoint.Registry where Registry.action=modified AND (Registry.registry_path=\"*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run*\" OR Registry.registry_path=\"*\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce*\") AND (Registry.registry_value_data=\"*\\\\AppData\\\\Roaming\\\\*\" OR Registry.registry_value_data=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Registry.registry_value_data=\"*\\\\Users\\\\Public\\\\*\") AND Registry.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\") by host Registry.dest Registry.user Registry.process_name Registry.process Registry.registry_path Registry.registry_value_name Registry.registry_value_data | `drop_dm_object_name(Registry)` | where NOT match(user, \"(?i).*\\\\$$\") | table firstSeen dest user process_name process registry_path registry_value_name registry_value_data",
      "kql": "DeviceRegistryEvents\n| where Timestamp > ago(7d)\n| where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has_any (@\"\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\",\n                            @\"\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\",\n                            @\"\\Software\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\")\n| where RegistryValueData has_any (@\"AppData\\Roaming\", @\"AppData\\Local\\Temp\", @\"\\Users\\Public\\\")\n| where InitiatingProcessFileName in~ (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\n                                       \"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\n                                       \"rundll32.exe\",\"regsvr32.exe\")\n| where InitiatingProcessAccountName !endswith \"$\"\n| project Timestamp, DeviceName, AccountName = InitiatingProcessAccountName,\n          RegistryKey, RegistryValueName, RegistryValueData,\n          ScriptInterpreter = InitiatingProcessFileName,\n          ScriptCmd = InitiatingProcessCommandLine,\n          GrandParent = InitiatingProcessParentFileName\n| order by Timestamp desc",
      "sentinel_kql": "ImRegistryEvent\n| where TimeGenerated > ago(7d)\n| where EventType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has_any (@\"\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\",\n                            @\"\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\",\n                            @\"\\Software\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\")\n| where RegistryValueData has_any (@\"AppData\\Roaming\", @\"AppData\\Local\\Temp\", @\"\\Users\\Public\\\")\n| where ActingProcessName has_any (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\n                                   \"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\n                                   \"rundll32.exe\",\"regsvr32.exe\")\n| where ActorUsername !endswith \"$\"\n| project TimeGenerated, DvcHostname, ActorUsername,\n          RegistryKey, RegistryValue, RegistryValueData,\n          ScriptInterpreter = ActingProcessName,\n          ScriptCmd = ActingProcessCommandLine\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Saint Bear / Saint Bot Run-Key Persistence via Script Interpreter\nid: 7c9d2f4a-1e3b-4b7d-9a6f-saintbear0002\nstatus: experimental\ndescription: Detects Saint Bear (UAC-0056 / Lorec53 / TA471) Saint Bot persistence \u2014 a script interpreter writes an HKCU Run/RunOnce value pointing at a binary in %AppData% or %Users%\\Public.\nreferences:\n  - https://attack.mitre.org/groups/G1031/\n  - https://unit42.paloaltonetworks.com/ukraine-targeted-outsteel-saintbot/\nauthor: SOC detection engineering\ndate: 2026/05/03\ntags:\n  - attack.persistence\n  - attack.t1547.001\n  - attack.t1112\n  - attack.g1031\nlogsource:\n  product: windows\n  category: registry_set\ndetection:\n  selection_key:\n    TargetObject|contains:\n      - '\\Software\\Microsoft\\Windows\\CurrentVersion\\Run'\n      - '\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce'\n      - '\\Software\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run'\n  selection_value:\n    Details|contains:\n      - '\\AppData\\Roaming\\'\n      - '\\AppData\\Local\\Temp\\'\n      - '\\Users\\Public\\'\n  selection_proc:\n    Image|endswith:\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n      - '\\mshta.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\cmd.exe'\n      - '\\rundll32.exe'\n      - '\\regsvr32.exe'\n  filter_machine:\n    User|endswith: '$'\n  condition: all of selection_* and not filter_machine\nfalsepositives:\n  - Portable applications that legitimately install themselves to %AppData% (rare; review parent chain)\nlevel: high",
      "rationale": "Saint Bot's persistence routine \u2014 well documented in Unit 42, Malwarebytes, and CERT-UA UAC-0056 reports \u2014 sets HKCU Run with a value pointing into %AppData% from a non-MSI, non-Setup parent (typically the same wscript/cscript that just dropped the binary in UC1). Legitimate software almost always writes Run keys via msiexec.exe, setup*.exe, or Explorer; a script interpreter doing it from a user c",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    }
  ]
}