# intel/.llm_actor_uc_cache/44/441c82123bfebba7320826f343d6dd9f5248054d.json
{
  "name": "Rocke",
  "ucs": [
    {
      "title": "[LLM] Rocke (Iron) \u2014 Linux cloud-agent uninstall preceding XMRig/Kinsing miner launch",
      "description": "Hunts Rocke's signature evasion pattern on Linux cloud workloads: the dropper uninstalls the Alibaba Cloud Aegis (/usr/local/aegis) and Tencent Cloud Monitor (/usr/local/qcloud, YunJing, barad_agent) agents before deploying its XMRig/cnrig/Kinsing-derived miner. This sequence \u2014 security-tooling removal directly followed by Monero mining on the same host within an hour \u2014 is publicly attributed to Rocke (Cisco Talos, 2018; Palo Alto Unit 42, 2019) and is rarely produced by legitimate admin activity.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1562.001",
        "T1518.001",
        "T1496.001",
        "T1059.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t min(_time) as agent_kill_time values(Processes.process) as agent_cmd from datamodel=Endpoint.Processes where Processes.os=Linux (Processes.process_name IN (bash,sh,systemctl,service,rm,pkill,kill,wget,curl) AND (Processes.process=\"*aegis_update*\" OR Processes.process=\"*aegis_cli*\" OR Processes.process=\"*/usr/local/aegis*\" OR Processes.process=\"*/usr/local/qcloud*\" OR Processes.process=\"*YunJing*\" OR Processes.process=\"*barad_agent*\" OR Processes.process=\"*CmsGoAgent*\" OR Processes.process=\"*quartz_update*\" OR Processes.process=\"*sgagent*\")) by Processes.dest Processes.user | rename Processes.* as * | join type=inner dest [| tstats summariesonly=t min(_time) as miner_time values(Processes.process) as miner_cmd values(Processes.process_name) as miner_proc from datamodel=Endpoint.Processes where Processes.os=Linux (Processes.process_name IN (xmrig,minerd,cnrig,kerberods,kthrotlds,sustes,kinsing,kdevtmpfsi) OR Processes.process=\"*stratum+tcp*\" OR Processes.process=\"*--donate-level*\" OR Processes.process=\"*cryptonight*\" OR Processes.process=\"*xmr.*\" OR Processes.process=\"*pool.minexmr*\" OR Processes.process=\"*supportxmr.com*\") by Processes.dest | rename Processes.* as *] | where miner_time>=agent_kill_time AND miner_time<=agent_kill_time+3600 | table dest user agent_kill_time agent_cmd miner_time miner_proc miner_cmd",
      "kql": "// Rocke (G0106) \u2014 cloud-agent removal followed by miner launch on the same Linux host within 1h\nlet Lookback = 14d;\nlet CloudAgentTokens = dynamic([\"aegis_update\",\"aegis_cli\",\"/usr/local/aegis\",\"/usr/local/qcloud\",\"YunJing\",\"yunjing\",\"barad_agent\",\"CmsGoAgent\",\"quartz_update\",\"sgagent\",\"AliYunDun\",\"AliHids\"]);\nlet MinerTokens = dynamic([\"xmrig\",\"minerd\",\"cnrig\",\"kerberods\",\"kthrotlds\",\"sustes\",\"kinsing\",\"kdevtmpfsi\",\"stratum+tcp\",\"--donate-level\",\"cryptonight\",\"pool.minexmr\",\"supportxmr.com\",\"xmr.\"]);\nlet AgentKills = DeviceProcessEvents\n    | where Timestamp > ago(Lookback)\n    | where FileName in~ (\"bash\",\"sh\",\"dash\",\"systemctl\",\"service\",\"rm\",\"pkill\",\"kill\",\"wget\",\"curl\")\n    | where ProcessCommandLine has_any (CloudAgentTokens)\n    | where ProcessCommandLine has_any (\"uninstall\",\"stop\",\"disable\",\" rm \",\"kill\",\"pkill\",\"systemctl stop\",\"service stop\")\n    | project AgentTime = Timestamp, DeviceId, DeviceName, AgentUser = AccountName,\n              AgentCmd = ProcessCommandLine, AgentParent = InitiatingProcessFileName,\n              AgentParentCmd = InitiatingProcessCommandLine;\nDeviceProcessEvents\n| where Timestamp > ago(Lookback)\n| where FileName has_any (MinerTokens) or ProcessCommandLine has_any (MinerTokens)\n| join kind=inner AgentKills on DeviceId\n| where Timestamp between (AgentTime .. AgentTime + 1h)\n| project DeviceName, AgentUser,\n          AgentTime, AgentCmd, AgentParent,\n          MinerTime = Timestamp, MinerProc = FileName, MinerCmd = ProcessCommandLine,\n          DelaySec = datetime_diff('second', Timestamp, AgentTime)\n| order by MinerTime desc",
      "sentinel_kql": "// Rocke (G0106) \u2014 Linux cloud-agent uninstall + miner launch within 1h, ASIM\nlet Lookback = 14d;\nlet CloudAgentTokens = dynamic([\"aegis_update\",\"aegis_cli\",\"/usr/local/aegis\",\"/usr/local/qcloud\",\"YunJing\",\"yunjing\",\"barad_agent\",\"CmsGoAgent\",\"quartz_update\",\"sgagent\",\"AliYunDun\",\"AliHids\"]);\nlet MinerTokens = dynamic([\"xmrig\",\"minerd\",\"cnrig\",\"kerberods\",\"kthrotlds\",\"sustes\",\"kinsing\",\"kdevtmpfsi\",\"stratum+tcp\",\"--donate-level\",\"cryptonight\",\"pool.minexmr\",\"supportxmr.com\",\"xmr.\"]);\nlet AgentKills = ImProcessCreate\n    | where TimeGenerated > ago(Lookback)\n    | where DvcOs has \"Linux\" or ActingProcessName in~ (\"bash\",\"sh\",\"dash\",\"systemctl\",\"service\")\n    | where TargetProcessName in~ (\"bash\",\"sh\",\"dash\",\"systemctl\",\"service\",\"rm\",\"pkill\",\"kill\",\"wget\",\"curl\")\n    | where TargetProcessCommandLine has_any (CloudAgentTokens)\n    | where TargetProcessCommandLine has_any (\"uninstall\",\"stop\",\"disable\",\" rm \",\"kill\",\"pkill\",\"systemctl stop\",\"service stop\")\n    | project AgentTime = TimeGenerated, DvcHostname, ActorUsername,\n              AgentCmd = TargetProcessCommandLine, AgentParent = ParentProcessName;\nImProcessCreate\n| where TimeGenerated > ago(Lookback)\n| where TargetProcessName has_any (MinerTokens) or TargetProcessCommandLine has_any (MinerTokens)\n| join kind=inner AgentKills on DvcHostname\n| where TimeGenerated between (AgentTime .. AgentTime + 1h)\n| project DvcHostname, ActorUsername,\n          AgentTime, AgentCmd, AgentParent,\n          MinerTime = TimeGenerated, MinerProc = TargetProcessName, MinerCmd = TargetProcessCommandLine,\n          DelaySec = datetime_diff('second', TimeGenerated, AgentTime)\n| order by MinerTime desc",
      "sigma_yaml": "title: Rocke Cloud Agent Uninstall (Aegis / Tencent qcloud)\nid: 7c9d4f6c-3b2a-4d18-bf7c-rocke-aegisuninstall\nstatus: experimental\ndescription: Detects Rocke (G0106) signature uninstall of Alibaba Cloud Aegis and Tencent Cloud Monitor agents on Linux, typically performed prior to XMRig deployment.\nreferences:\n  - https://blog.talosintelligence.com/rocke-champion-of-monero-miners/\n  - https://attack.mitre.org/groups/G0106/\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.defense_evasion\n  - attack.t1562.001\n  - attack.t1518.001\n  - attack.impact\n  - attack.t1496.001\nlogsource:\n  product: linux\n  category: process_creation\ndetection:\n  selection_proc:\n    Image|endswith:\n      - '/bash'\n      - '/sh'\n      - '/dash'\n      - '/systemctl'\n      - '/service'\n      - '/rm'\n      - '/pkill'\n      - '/kill'\n      - '/wget'\n      - '/curl'\n  selection_target:\n    CommandLine|contains:\n      - 'aegis_update'\n      - 'aegis_cli'\n      - '/usr/local/aegis'\n      - '/usr/local/qcloud'\n      - 'YunJing'\n      - 'barad_agent'\n      - 'CmsGoAgent'\n      - 'quartz_update'\n      - 'AliYunDun'\n      - 'AliHids'\n  selection_action:\n    CommandLine|contains:\n      - 'uninstall'\n      - 'stop'\n      - 'disable'\n      - ' rm '\n      - 'kill'\n      - 'pkill'\n  condition: selection_proc and selection_target and selection_action\nfalsepositives:\n  - Authorised cloud-agent migration / re-imaging (rare; should be ticketed)\nlevel: high",
      "rationale": "Talos and Unit 42 specifically documented Rocke removing Aegis and qcloud agents in its first-stage shell script (a0.sh / 2t3ik / xmr_cgmin variants) immediately before unpacking the miner. Random opportunistic cryptojackers usually skip cloud-agent removal \u2014 only Rocke and the closely-related Iron Cybercrime Group ship the Alibaba/Tencent uninstall logic. Pairing the agent-kill with a miner spawn",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] Rocke \u2014 Linux web-app exploit (Struts/WebLogic/ColdFusion) spawning shell that pulls dropper from Pastebin/Gitee with cron persistence",
      "description": "Hunts Rocke's full initial-access chain: Apache Struts / Oracle WebLogic (CVE-2017-10271) / Adobe ColdFusion (CVE-2017-3066) RCE \u2192 java/php-fpm/httpd spawns bash \u2192 curl/wget to pastebin.com, gitee.com, transfer.sh or pastes-dot-ee \u2192 chmod +x \u2192 crontab append for /etc/cron.d, /var/spool/cron, or .bashrc tail. Rocke characteristically uses Pastebin/Gitee as a dead-drop resolver for first-stage shell scripts (T1102.001) and installs a cron entry within the same one-liner.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1190",
        "T1059.004",
        "T1105",
        "T1102.001",
        "T1053.003",
        "T1222.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as first_seen values(Processes.process) as cmdlines values(Processes.parent_process_name) as parents values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.os=Linux Processes.parent_process_name IN (java,javaw,httpd,apache2,nginx,php-fpm,php,tomcat,weblogic,coldfusion,node,python,python3,ruby) Processes.process_name IN (bash,sh,dash,wget,curl,perl,python,python3) (Processes.process=\"*pastebin.com*\" OR Processes.process=\"*gitee.com*\" OR Processes.process=\"*raw.githubusercontent*\" OR Processes.process=\"*transfer.sh*\" OR Processes.process=\"*gitlab.com/*/raw*\" OR Processes.process=\"*paste.ee*\") (Processes.process=\"*chmod *\" OR Processes.process=\"*crontab*\" OR Processes.process=\"*/etc/cron*\" OR Processes.process=\"*/var/spool/cron*\" OR Processes.process=\"*.bashrc*\" OR Processes.process=\"* | sh*\" OR Processes.process=\"* | bash*\") by Processes.dest Processes.process_name | rename Processes.* as * | sort 0 - first_seen",
      "kql": "// Rocke (G0106) \u2014 public-facing app exploit drops loader from Pastebin/Gitee with cron persistence one-liner\nlet Lookback = 7d;\nlet WebParents = dynamic([\"java\",\"javaw\",\"httpd\",\"apache2\",\"nginx\",\"php-fpm\",\"php\",\"tomcat\",\"weblogic\",\"coldfusion\",\"node\",\"python\",\"python3\",\"ruby\"]);\nlet DeadDrops = dynamic([\"pastebin.com\",\"gitee.com\",\"raw.githubusercontent.com\",\"transfer.sh\",\"paste.ee\",\"gitlab.com\",\"hastebin.com\",\"sprunge.us\",\"thyrsi.com\"]);\nlet PersistTokens = dynamic([\"chmod +x\",\"chmod 777\",\"chmod 755\",\"crontab\",\"/etc/cron\",\"/var/spool/cron\",\".bashrc\",\".profile\",\"/etc/rc.local\",\" | sh\",\" | bash\",\"systemctl enable\"]);\nDeviceProcessEvents\n| where Timestamp > ago(Lookback)\n| where InitiatingProcessFileName in~ (WebParents)\n| where FileName in~ (\"bash\",\"sh\",\"dash\",\"wget\",\"curl\",\"perl\",\"python\",\"python3\")\n| where ProcessCommandLine has_any (DeadDrops)\n| where ProcessCommandLine has_any (PersistTokens)\n| project Timestamp, DeviceName, AccountName,\n          WebParent = InitiatingProcessFileName,\n          WebParentCmd = InitiatingProcessCommandLine,\n          LoaderProc = FileName,\n          LoaderCmd = ProcessCommandLine,\n          SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// Rocke (G0106) \u2014 Struts/WebLogic/ColdFusion exploit drops Pastebin/Gitee loader + cron, ASIM Linux\nlet Lookback = 7d;\nlet WebParents = dynamic([\"java\",\"javaw\",\"httpd\",\"apache2\",\"nginx\",\"php-fpm\",\"php\",\"tomcat\",\"weblogic\",\"coldfusion\",\"node\",\"python\",\"python3\",\"ruby\"]);\nlet DeadDrops = dynamic([\"pastebin.com\",\"gitee.com\",\"raw.githubusercontent.com\",\"transfer.sh\",\"paste.ee\",\"gitlab.com\",\"hastebin.com\",\"sprunge.us\",\"thyrsi.com\"]);\nlet PersistTokens = dynamic([\"chmod +x\",\"chmod 777\",\"chmod 755\",\"crontab\",\"/etc/cron\",\"/var/spool/cron\",\".bashrc\",\".profile\",\"/etc/rc.local\",\" | sh\",\" | bash\",\"systemctl enable\"]);\nImProcessCreate\n| where TimeGenerated > ago(Lookback)\n| where ParentProcessName in~ (WebParents)\n| where TargetProcessName in~ (\"bash\",\"sh\",\"dash\",\"wget\",\"curl\",\"perl\",\"python\",\"python3\")\n| where TargetProcessCommandLine has_any (DeadDrops)\n| where TargetProcessCommandLine has_any (PersistTokens)\n| project TimeGenerated, DvcHostname, ActorUsername,\n          WebParent = ParentProcessName,\n          WebParentCmd = ParentProcessCommandLine,\n          LoaderProc = TargetProcessName,\n          LoaderCmd = TargetProcessCommandLine,\n          SHA256\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "Talos's 2018 'Rocke: The Champion of Monero Miners' report and Anomali's 2019 follow-up document Rocke's invariant initial-access chain almost exactly: a vulnerable Java/PHP/CF web service is exploited, the resulting shell uses curl/wget against pastebin.com or gitee.com (later switched to gitee after Pastebin takedowns of its accounts) to fetch a0.sh / 1.sh / 2t3ik, immediately chmods it, and app"
    }
  ]
}