# intel/.llm_actor_uc_cache/44/44954ac375133d64b4d8c36c1dd9891957832506.json
{
  "name": "PROMETHIUM",
  "ucs": [
    {
      "title": "[LLM] PROMETHIUM/StrongPity trojanized installer dropper chain (Downloads\u2192%AppData% drop\u2192Run/Service persistence)",
      "description": "Hunts the signature StrongPity delivery pattern: a user-executed installer masquerading as a legitimate utility (WinRAR, 7-Zip, TrueCrypt, VPN client, Notepad++, etc.) \u2014 typically delivered via watering-hole / lookalike download portals or ISP-level redirection \u2014 that drops a secondary executable into %AppData%/%ProgramData% within minutes and immediately establishes Run-key or service persistence. This three-stage temporal chain is characteristic of PROMETHIUM/StrongPity droppers and is rarely seen from benign installers, which write to %ProgramFiles% and use the legitimate Windows Installer ",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1189",
        "T1204.002",
        "T1036.005",
        "T1553.002",
        "T1587.002",
        "T1543.003",
        "T1547.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as installerTime values(Processes.process) as installer_cmd values(Processes.process_path) as installer_path from datamodel=Endpoint.Processes where (Processes.parent_process_path IN (\"*\\\\Downloads\\\\*\",\"*\\\\Temp\\\\*\",\"*\\\\AppData\\\\Local\\\\Temp\\\\*\") OR Processes.process_path IN (\"*\\\\Downloads\\\\*\",\"*\\\\Temp\\\\*\")) AND (Processes.process_name IN (\"*setup*.exe\",\"*install*.exe\",\"wrar*.exe\",\"winrar*.exe\",\"7z*.exe\",\"truecrypt*.exe\",\"veracrypt*.exe\",\"*vpn*.exe\",\"npp*.exe\",\"notepad++*.exe\",\"firefox*.exe\",\"chrome*setup*.exe\",\"teamviewer*.exe\",\"skype*.exe\",\"driverpack*.exe\",\"ccleaner*.exe\",\"opera*.exe\")) by host Processes.user Processes.process_name Processes.process_guid | `drop_dm_object_name(Processes)` | rename process_guid as installer_guid | join type=inner host user [ | tstats summariesonly=true count min(_time) as dropTime values(Filesystem.file_path) as dropped_path values(Filesystem.file_name) as dropped_file from datamodel=Endpoint.Filesystem where Filesystem.file_path IN (\"*\\\\AppData\\\\Roaming\\\\*\",\"*\\\\AppData\\\\Local\\\\*\",\"*\\\\ProgramData\\\\*\") AND (Filesystem.file_name=\"*.exe\" OR Filesystem.file_name=\"*.dll\") by host Filesystem.user Filesystem.process_guid | `drop_dm_object_name(Filesystem)` | rename process_guid as installer_guid ] | join type=inner host user [ | tstats summariesonly=true count min(_time) as persistTime values(Registry.registry_path) as persist_key values(Registry.registry_value_data) as persist_value from datamodel=Endpoint.Registry where Registry.registry_path IN (\"*\\\\CurrentVersion\\\\Run\\\\*\",\"*\\\\CurrentVersion\\\\RunOnce\\\\*\",\"*\\\\CurrentControlSet\\\\Services\\\\*\\\\ImagePath\") by host Registry.user ] | eval drop_lag=dropTime-installerTime, persist_lag=persistTime-installerTime | where drop_lag>=0 AND drop_lag<=600 AND persist_lag>=0 AND persist_lag<=1800 | eval risk_signal=\"PROMETHIUM_StrongPity_dropper_chain\" | table installerTime host user installer_cmd installer_path dropped_file dropped_path persist_key persist_value drop_lag persist_lag",
      "kql": "let installer_tokens = dynamic([\"setup\",\"install\",\"winrar\",\"wrar\",\"7zip\",\"7z\",\"truecrypt\",\"veracrypt\",\"vpn\",\"npp\",\"notepad++\",\"firefox\",\"chrome\",\"teamviewer\",\"skype\",\"driverpack\",\"ccleaner\",\"opera\"]);\nlet droppers =\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where FolderPath has_any (@\"\\Downloads\\\", @\"\\Temp\\\", @\"\\AppData\\Local\\Temp\\\")\n    | where FileName has_any (installer_tokens)\n    | project InstallerTime=Timestamp, DeviceId, DeviceName, AccountName=InitiatingProcessAccountName,\n              InstallerName=FileName, InstallerPath=FolderPath, InstallerSha=SHA256,\n              InstallerProcessId=ProcessId, InstallerCmd=ProcessCommandLine;\nlet drops =\n    DeviceFileEvents\n    | where Timestamp > ago(7d)\n    | where ActionType == \"FileCreated\"\n    | where FolderPath has_any (@\"\\AppData\\Roaming\\\", @\"\\AppData\\Local\\\", @\"\\ProgramData\\\")\n    | where FileName endswith \".exe\" or FileName endswith \".dll\"\n    | where not(FolderPath has_any (@\"\\Microsoft\\\", @\"\\Packages\\\", @\"\\WindowsApps\\\"))\n    | project DropTime=Timestamp, DeviceId, DroppedFile=FileName, DroppedPath=FolderPath,\n              DroppedSha=SHA256, InitiatingProcessId, InitiatingProcessFileName;\nlet persistence =\n    union\n    (DeviceRegistryEvents\n        | where Timestamp > ago(7d)\n        | where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n        | where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\", @\"\\Services\\\")\n        | project PersistTime=Timestamp, DeviceId, PersistArtifact=strcat(RegistryKey,\"=\",RegistryValueData), PersistInitiator=InitiatingProcessId),\n    (DeviceProcessEvents\n        | where Timestamp > ago(7d)\n        | where FileName =~ \"sc.exe\" and ProcessCommandLine has \"create\"\n        | project PersistTime=Timestamp, DeviceId, PersistArtifact=ProcessCommandLine, PersistInitiator=InitiatingProcessId);\ndroppers\n| join kind=inner drops on DeviceId\n| where DropTime between (InstallerTime .. InstallerTime + 10m)\n| where InitiatingProcessId == InstallerProcessId or InitiatingProcessFileName =~ InstallerName\n| join kind=inner persistence on DeviceId\n| where PersistTime between (InstallerTime .. InstallerTime + 30m)\n| extend DropLagSec = datetime_diff('second', DropTime, InstallerTime),\n         PersistLagSec = datetime_diff('second', PersistTime, InstallerTime),\n         ActorHypothesis = \"PROMETHIUM / StrongPity (G0056) dropper chain\"\n| project InstallerTime, DeviceName, AccountName, InstallerName, InstallerPath, InstallerSha,\n          DroppedFile, DroppedPath, DroppedSha, PersistArtifact, DropLagSec, PersistLagSec, ActorHypothesis\n| sort by InstallerTime desc",
      "rationale": "This isn't a generic 'installer drops file' rule \u2014 it encodes the exact StrongPity tradecraft Talos and Microsoft (MSTIC) documented in the 2016 SIR Vol 21 NEODYMIUM/PROMETHIUM report and the June-2020 Talos StrongPity4 writeup: PROMETHIUM clones publisher download portals (and at times leverages ISP-level injection on T\u00fcrk Telekom infrastructure), serves a working but trojanized version of common"
    }
  ]
}