# intel/.llm_actor_uc_cache/42/42a2cc8ef8f6d46d2f95ec2d929194805aeac68b.json
{
  "name": "APT41",
  "ucs": [
    {
      "title": "[LLM] APT41 (Wicked Panda) DLL side-loading via known signed lure binary from staging path",
      "description": "APT41/Winnti's signature defense-evasion is dropping a legitimate signed third-party binary (DTHelper.exe from Mediatek BT, oleview.exe, BdReinit.exe, sllauncher.exe, chrome_frame_helper.exe, log.exe, mcods.exe, etc.) into a writable staging directory and side-loading a malicious DLL (ShadowPad / PlugX / Winnti core / KEYPLUG / SPYDER) from the same folder. This UC catches that exact tradecraft: a known APT41 lure binary running from a non-vendor path, with a sibling unsigned DLL imported within seconds.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.002",
        "T1036.005",
        "T1027.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) AS firstTime max(_time) AS lastTime FROM datamodel=Endpoint.Processes WHERE Processes.process_name IN (\"BdReinit.exe\",\"oleview.exe\",\"DTHelper.exe\",\"sllauncher.exe\",\"chrome_frame_helper.exe\",\"HPCustParticUI.exe\",\"log.exe\",\"mcods.exe\",\"vmnat.exe\") AND (Processes.process_path LIKE \"%\\\\ProgramData\\\\%\" OR Processes.process_path LIKE \"%\\\\Temp\\\\%\" OR Processes.process_path LIKE \"%\\\\AppData\\\\%\" OR Processes.process_path LIKE \"%\\\\Public\\\\%\" OR Processes.process_path LIKE \"%\\\\PerfLogs\\\\%\" OR Processes.process_path LIKE \"%\\\\Intel\\\\%\") AND NOT match(Processes.user,\"\\$$\") BY host Processes.user Processes.process_name Processes.process_path Processes.process Processes.process_hash Processes.parent_process_name Processes.dest | `drop_dm_object_name(Processes)` | join type=inner host process_path [ | tstats summariesonly=true count AS dll_drop_count, min(_time) AS dll_time FROM datamodel=Endpoint.Filesystem WHERE Filesystem.file_name=\"*.dll\" AND (Filesystem.file_path LIKE \"%\\\\ProgramData\\\\%\" OR Filesystem.file_path LIKE \"%\\\\Temp\\\\%\" OR Filesystem.file_path LIKE \"%\\\\AppData\\\\%\" OR Filesystem.file_path LIKE \"%\\\\Public\\\\%\" OR Filesystem.file_path LIKE \"%\\\\PerfLogs\\\\%\") BY host Filesystem.file_name Filesystem.file_path Filesystem.file_hash | `drop_dm_object_name(Filesystem)` | rename file_path AS process_path file_name AS sideloaded_dll file_hash AS sideloaded_dll_sha256 ] | where dll_time <= firstTime + 600 AND dll_time >= firstTime - 1800 | convert ctime(firstTime) ctime(dll_time) | table firstTime, host, user, process_name, process_path, process_hash, parent_process_name, process, sideloaded_dll, sideloaded_dll_sha256",
      "kql": "let LookbackDays = 7d;\nlet _apt41_loaders = dynamic([\"bdreinit.exe\",\"oleview.exe\",\"dthelper.exe\",\"sllauncher.exe\",\"chrome_frame_helper.exe\",\"hpcustparticui.exe\",\"log.exe\",\"mcods.exe\",\"vmnat.exe\"]);\nlet _staging_paths = dynamic([@\"\\ProgramData\\\", @\"\\Temp\\\", @\"\\AppData\\\", @\"\\Public\\\", @\"\\PerfLogs\\\", @\"\\Intel\\\"]);\nlet Loaders = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where tolower(FileName) in (_apt41_loaders)\n    | where FolderPath has_any (_staging_paths)\n    | where FolderPath !startswith @\"C:\\Program Files\" and FolderPath !startswith @\"C:\\Program Files (x86)\"\n    | project LoaderTime=Timestamp, DeviceId, DeviceName, AccountName,\n              LoaderName=FileName, LoaderPath=FolderPath, LoaderSHA256=SHA256,\n              LoaderCmd=ProcessCommandLine,\n              ParentName=InitiatingProcessFileName, ParentCmd=InitiatingProcessCommandLine;\nDeviceImageLoadEvents\n| where Timestamp > ago(LookbackDays)\n| where FileName endswith \".dll\"\n| where FolderPath !startswith @\"C:\\Windows\\\"\n| where FolderPath has_any (_staging_paths)\n| join kind=inner Loaders on DeviceId, $left.InitiatingProcessFileName == $right.LoaderName\n| where Timestamp between (LoaderTime .. LoaderTime + 5m)\n| where FolderPath =~ LoaderPath   // sibling DLL = side-load\n| where SHA256 != LoaderSHA256\n| project LoaderTime, ImageLoadTime=Timestamp, DeviceName, AccountName,\n          LoaderName, LoaderPath, LoaderSHA256,\n          SideloadedDll=FileName, SideloadedDllSHA256=SHA256,\n          ParentName, ParentCmd, LoaderCmd\n| order by LoaderTime desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet _apt41_loaders = dynamic([\"BdReinit.exe\",\"oleview.exe\",\"DTHelper.exe\",\"sllauncher.exe\",\"chrome_frame_helper.exe\",\"HPCustParticUI.exe\",\"log.exe\",\"mcods.exe\",\"vmnat.exe\"]);\nlet _staging_paths = dynamic([@\"\\ProgramData\\\", @\"\\Temp\\\", @\"\\AppData\\\", @\"\\Public\\\", @\"\\PerfLogs\\\", @\"\\Intel\\\"]);\nlet Loaders = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend LoaderShort = tostring(split(NewProcessName, @\"\\\")[-1])\n    | where LoaderShort in~ (_apt41_loaders)\n    | where NewProcessName has_any (_staging_paths)\n    | where NewProcessName !startswith @\"C:\\Program Files\"\n    | project LoaderTime=TimeGenerated, Computer, Account, NewProcessName,\n              LoaderShort, CommandLine, ParentProcessName;\n// Pivot to Sysmon image-load (EID 7) for the actual side-loaded DLL when available\nlet ImageLoads = WindowsEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where Provider == \"Microsoft-Windows-Sysmon\" and EventID == 7\n    | extend Image = tostring(EventData.Image),\n             ImageLoaded = tostring(EventData.ImageLoaded),\n             Signed = tostring(EventData.Signed),\n             Hashes = tostring(EventData.Hashes)\n    | extend LoaderShort = tostring(split(Image, @\"\\\")[-1])\n    | where LoaderShort in~ (_apt41_loaders)\n    | where ImageLoaded endswith \".dll\"\n    | where ImageLoaded has_any (_staging_paths);\nLoaders\n| join kind=inner ImageLoads on Computer, LoaderShort\n| where TimeGenerated between (LoaderTime .. LoaderTime + 5m)\n| project LoaderTime, ImageLoadTime=TimeGenerated, Computer, Account,\n          Loader=NewProcessName, LoaderCmd=CommandLine, ParentProcessName,\n          SideloadedDll=ImageLoaded, Signed, Hashes\n| order by LoaderTime desc",
      "sigma_yaml": "title: APT41 DLL Side-Loading via Known Lure Binary from Staging Path\nid: 8b3a4d20-a141-4e5f-bf41-9b7c12a5e001\nstatus: experimental\ndescription: Detects execution of legitimate third-party signed binaries that APT41 (Wicked Panda / BARIUM / Winnti) repeatedly abuses for DLL side-loading of ShadowPad, PlugX, Winnti and KEYPLUG, when launched from non-vendor staging directories.\nreferences:\n  - https://www.mandiant.com/resources/blog/apt41-us-state-governments\n  - https://attack.mitre.org/groups/G0096/\nauthor: SOC Detection Engineering\ndate: 2026-05-07\ntags:\n  - attack.defense_evasion\n  - attack.t1574.002\n  - attack.t1036.005\n  - attack.g0096\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_lure:\n    Image|endswith:\n      - '\\BdReinit.exe'\n      - '\\oleview.exe'\n      - '\\DTHelper.exe'\n      - '\\sllauncher.exe'\n      - '\\chrome_frame_helper.exe'\n      - '\\HPCustParticUI.exe'\n      - '\\log.exe'\n      - '\\mcods.exe'\n      - '\\vmnat.exe'\n  selection_path:\n    Image|contains:\n      - '\\ProgramData\\'\n      - '\\Temp\\'\n      - '\\AppData\\'\n      - '\\Public\\'\n      - '\\PerfLogs\\'\n      - '\\Intel\\'\n  filter_legit:\n    Image|startswith:\n      - 'C:\\Program Files\\'\n      - 'C:\\Program Files (x86)\\'\n  condition: selection_lure and selection_path and not filter_legit\nfalsepositives:\n  - Rare: legitimate copies of these vendor utilities staged outside install directories by IT for portable use\nlevel: high",
      "rationale": "APT41 has been pinned to side-loading via this exact short-list of carrier binaries across Mandiant, Kaspersky and ESET reporting (DTHelper.exe\u2192ShadowPad, sllauncher.exe\u2192PlugX, oleview.exe & BdReinit.exe\u2192Winnti, log.exe\u2192KEYPLUG). The combination of (a) one of these specific binary names, (b) running outside its normal vendor install path, and (c) loading a sibling DLL within seconds is essentially"
    },
    {
      "title": "[LLM] APT41 web shell on Exchange/IIS chained to scheduled-task persistence (China Chopper / Antsword to schtasks)",
      "description": "APT41/BARIUM consistently exploits public-facing IIS or Exchange (e.g. ProxyLogon, Citrix, Log4Shell) to drop China Chopper / Antsword / regeorg web shells, then immediately runs recon (whoami, net, systeminfo) and creates a scheduled task as persistence (often masquerading under Microsoft\\Windows\\* paths). This UC fires when w3wp.exe spawns recon/LOLBin children AND a schtasks /create executes on the same host within 30 minutes \u2014 the hands-on-keyboard footprint that distinguishes APT41 from a noisy commodity web-shell scan.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1505.003",
        "T1190",
        "T1053.005",
        "T1059.001",
        "T1059.003",
        "T1036.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) AS shell_time max(_time) AS shell_last FROM datamodel=Endpoint.Processes WHERE Processes.parent_process_name IN (\"w3wp.exe\",\"httpd.exe\",\"tomcat.exe\",\"java.exe\") AND Processes.process_name IN (\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"ipconfig.exe\",\"nltest.exe\",\"quser.exe\",\"reg.exe\",\"sc.exe\",\"wmic.exe\",\"schtasks.exe\",\"certutil.exe\",\"bitsadmin.exe\") AND NOT match(Processes.user,\"\\$$\") BY host Processes.user Processes.parent_process_name Processes.process_name Processes.process Processes.process_path | `drop_dm_object_name(Processes)` | join type=inner host [ | tstats summariesonly=true count AS sch_count min(_time) AS sch_time values(Processes.process) AS sch_cmd values(Processes.parent_process_name) AS sch_parent FROM datamodel=Endpoint.Processes WHERE Processes.process_name=\"schtasks.exe\" AND (Processes.process LIKE \"%/create%\" OR Processes.process LIKE \"%/CREATE%\" OR Processes.process LIKE \"%-Create%\") BY host ] | where sch_time >= shell_time AND sch_time <= shell_time + 1800 | eval delay_sec = sch_time - shell_time | convert ctime(shell_time) ctime(sch_time) | table shell_time, host, user, parent_process_name, process_name, process, sch_time, delay_sec, sch_parent, sch_cmd",
      "kql": "let LookbackDays = 7d;\nlet _shell_tools = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"ipconfig.exe\",\"nltest.exe\",\"quser.exe\",\"reg.exe\",\"sc.exe\",\"wmic.exe\",\"schtasks.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"hostname.exe\"]);\nlet _web_servers = dynamic([\"w3wp.exe\",\"httpd.exe\",\"tomcat.exe\",\"java.exe\",\"nginx.exe\"]);\nlet WebShellChildren = DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where InitiatingProcessFileName in~ (_web_servers)\n    | where FileName in~ (_shell_tools)\n    // Strip baseline Exchange/IIS health-probe noise\n    | where not(InitiatingProcessCommandLine has_any (\"MSExchangeHealthManager\",\"MSExchangeOWAAppPool\") and ProcessCommandLine has \"DefaultAppPool\")\n    | project ShellTime=Timestamp, DeviceId, DeviceName, AccountName,\n              IISParent=InitiatingProcessFileName, IISCmd=InitiatingProcessCommandLine,\n              ChildProc=FileName, ChildCmd=ProcessCommandLine;\nDeviceProcessEvents\n| where Timestamp > ago(LookbackDays)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has_any (\"/create\",\"/Create\",\"/CREATE\")\n| project SchTaskTime=Timestamp, DeviceId, DeviceName,\n          SchTaskCmd=ProcessCommandLine,\n          SchInitiator=InitiatingProcessFileName, SchInitiatorCmd=InitiatingProcessCommandLine\n| join kind=inner WebShellChildren on DeviceId\n| where SchTaskTime between (ShellTime .. ShellTime + 30m)\n| extend DelaySec = datetime_diff('second', SchTaskTime, ShellTime)\n| extend LooksMasqueraded = SchTaskCmd has_any (@\"\\Microsoft\\\\Windows\\\\\", \"WindowsUpdate\", \"DefenderUpdate\", \"OneDriveSetup\", \"AdobeFlashPlayer\")\n| project ShellTime, DeviceName, AccountName, IISParent, IISCmd,\n          ChildProc, ChildCmd, SchTaskTime, DelaySec, SchTaskCmd, SchInitiator,\n          LooksMasqueraded\n| order by ShellTime desc",
      "sentinel_kql": "let LookbackDays = 7d;\nlet _shell_tools = dynamic([\"cmd.exe\",\"powershell.exe\",\"pwsh.exe\",\"cscript.exe\",\"wscript.exe\",\"whoami.exe\",\"net.exe\",\"net1.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"ipconfig.exe\",\"nltest.exe\",\"quser.exe\",\"reg.exe\",\"sc.exe\",\"wmic.exe\",\"schtasks.exe\",\"certutil.exe\",\"bitsadmin.exe\"]);\nlet _web_servers = dynamic([\"w3wp.exe\",\"httpd.exe\",\"tomcat.exe\",\"java.exe\",\"nginx.exe\"]);\nlet WebShellChildren = SecurityEvent\n    | where TimeGenerated > ago(LookbackDays)\n    | where EventID == 4688\n    | where Account !endswith \"$\"\n    | extend ParentShort = tostring(split(ParentProcessName, @\"\\\")[-1]),\n             ChildShort  = tostring(split(NewProcessName,    @\"\\\")[-1])\n    | where ParentShort in~ (_web_servers)\n    | where ChildShort in~ (_shell_tools)\n    | project ShellTime=TimeGenerated, Computer, Account,\n              ParentShort, ChildShort, NewProcessName, CommandLine, ParentProcessName;\nSecurityEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventID == 4688\n| extend ChildShort = tostring(split(NewProcessName, @\"\\\")[-1])\n| where ChildShort =~ \"schtasks.exe\"\n| where CommandLine has_any (\"/create\",\"/Create\",\"/CREATE\")\n| project SchTaskTime=TimeGenerated, Computer, SchTaskCmd=CommandLine,\n          SchParent=ParentProcessName\n| join kind=inner WebShellChildren on Computer\n| where SchTaskTime between (ShellTime .. ShellTime + 30m)\n| extend DelaySec = datetime_diff('second', SchTaskTime, ShellTime)\n| extend LooksMasqueraded = SchTaskCmd has_any (@\"\\Microsoft\\Windows\\\", \"WindowsUpdate\", \"DefenderUpdate\", \"OneDriveSetup\", \"AdobeFlashPlayer\")\n| project ShellTime, Computer, Account, ParentShort, ChildShort, CommandLine,\n          SchTaskTime, DelaySec, SchTaskCmd, SchParent, LooksMasqueraded\n| order by ShellTime desc",
      "sigma_yaml": "",
      "rationale": "APT41 is publicly attributed to multiple ProxyLogon (CVE-2021-26855) and follow-on Exchange/IIS web-shell campaigns (Mandiant: 'APT41 Targets U.S. State Governments', 2022; CISA AA22-138B). Their post-exploitation playbook is consistently: w3wp child process \u2192 recon (whoami, systeminfo, net group) \u2192 schtasks /create with a Microsoft-style task path for ShadowPad/HIGHNOON persistence. Joining a w3w"
    }
  ]
}