# intel/.llm_actor_uc_cache/78/784693db970bf0be8a6b0786753da10e40db7f3a.json
{
  "name": "Higaisa",
  "ucs": [
    {
      "title": "[LLM] Higaisa LNK-borne spearphish: cmd.exe carving embedded payload from .lnk via findstr/certutil with decoy PDF drop",
      "description": "Hunts Higaisa's signature initial-access chain in which a spearphished ZIP contains a double-extension .pdf.lnk. When opened, the LNK invokes cmd.exe to self-extract an embedded base64/PE blob from the same .lnk file using findstr/certutil, then opens a decoy PDF to mask execution. This pattern was the core delivery TTP across the 2020 Malwarebytes/Zscaler/PT Security disclosures of G0126.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.003",
        "T1140",
        "T1564.003",
        "T1036.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as process values(Processes.process_name) as process_name values(Processes.parent_process) as parent_cmdline values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"explorer.exe\",\"OUTLOOK.EXE\",\"WINWORD.EXE\") AND Processes.process_name=\"cmd.exe\" AND (Processes.process=\"*findstr*\\\".lnk*\" OR Processes.process=\"*for /f*.lnk*\" OR (Processes.process=\"*certutil*-decode*\" AND Processes.process=\"*.lnk*\") OR Processes.process=\"*start /min*\") by Processes.dest Processes.parent_process_id Processes.process_id _time | `drop_dm_object_name(Processes)` | join type=inner dest [| tstats summariesonly=t count from datamodel=Endpoint.Filesystem where Filesystem.file_name=\"*.pdf\" AND Filesystem.file_path IN (\"*\\\\Temp\\\\*\",\"*\\\\AppData\\\\Local\\\\Temp\\\\*\",\"*\\\\Public\\\\*\") by Filesystem.dest Filesystem.file_name Filesystem.file_path | `drop_dm_object_name(Filesystem)` | rename file_name as decoy_pdf] | where lastTime - firstTime < 120 | table firstTime dest user parent_cmdline process decoy_pdf",
      "kql": "let lnkCarvers = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (\"explorer.exe\",\"outlook.exe\",\"winword.exe\",\"7zfm.exe\",\"winrar.exe\")\n| where FileName =~ \"cmd.exe\"\n| where ProcessCommandLine has \".lnk\"\n    and (ProcessCommandLine has_any (\"findstr\",\"for /f\",\"certutil -decode\",\"certutil /decode\")\n         or ProcessCommandLine has_all (\"start\",\"/min\"))\n| project carveTime=Timestamp, DeviceId, DeviceName, AccountName, lnkCmd=ProcessCommandLine, ParentImg=InitiatingProcessFileName, ParentCmd=InitiatingProcessCommandLine;\nlet decoyDrops = DeviceFileEvents\n| where Timestamp > ago(7d)\n| where ActionType == \"FileCreated\"\n| where FileName endswith \".pdf\"\n| where FolderPath has_any (@\"\\Temp\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Public\\\")\n| project dropTime=Timestamp, DeviceId, decoyPath=FolderPath, decoyName=FileName, droppingProcess=InitiatingProcessFileName;\nlnkCarvers\n| join kind=inner decoyDrops on DeviceId\n| where (dropTime - carveTime) between (-2m .. 2m)\n| project carveTime, DeviceName, AccountName, ParentImg, lnkCmd, decoyName, decoyPath",
      "rationale": "Generic LNK-execution alerts are noisy. Higaisa's actual operational pattern is uniquely narrow: a user-launched LNK spawns cmd.exe which references the LNK file by name via findstr/certutil to carve an embedded payload, with a near-simultaneous PDF decoy drop in Temp/Public. Requiring the LNK-as-data-source string AND a decoy PDF write within ~120s on the same host targets G0126's documented deli"
    },
    {
      "title": "[LLM] Higaisa scheduled-task persistence: schtasks creating recurring jobs that launch wscript/mshta with hidden windows from %AppData%/%Public%",
      "description": "Hunts Higaisa's post-exploitation persistence: after the LNK detonates, the actor establishes a scheduled task that re-launches a JavaScript/VBScript or HTA stager from a user-writable path with /min or hidden-window flags. This combo of T1053.005 + T1059.005/.007 + T1564.003 + T1220 targeting %AppData%/%Public% paths is consistently observed across G0126 intrusions.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "medium",
      "techs": [
        "T1053.005",
        "T1059.005",
        "T1059.007",
        "T1220",
        "T1564.003",
        "T1547.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmd values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" AND Processes.process=\"*/create*\" AND (Processes.process=\"*wscript*\" OR Processes.process=\"*cscript*\" OR Processes.process=\"*mshta*\" OR Processes.process=\"*wmic*os*get*/format*\") AND (Processes.process=\"*\\\\AppData\\\\*\" OR Processes.process=\"*\\\\Public\\\\*\" OR Processes.process=\"*\\\\ProgramData\\\\*\" OR Processes.process=\"*\\\\Users\\\\Public\\\\*\") by Processes.dest Processes.user Processes.parent_process_name | `drop_dm_object_name(Processes)` | rex field=cmd \"/tn\\s+\\\"?(?<task_name>[^\\\"\\s]+)\" | rex field=cmd \"/tr\\s+\\\"?(?<task_run>[^\\\"]+)\" | where NOT match(parent,\"(?i)(SCCM|ccmexec|gpscript|TiWorker)\\.exe\") | table firstTime dest user parent task_name task_run cmd",
      "kql": "DeviceProcessEvents\n| where Timestamp > ago(14d)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has_cs \"/create\"\n| where ProcessCommandLine has_any (\"wscript\",\"cscript\",\"mshta\",\"wmic os get /format\")\n| where ProcessCommandLine has_any (@\"\\AppData\\\", @\"\\Public\\\", @\"\\ProgramData\\\", @\"\\Users\\Public\\\")\n| extend taskName = extract(@\"/tn\\s+\\\"?([^\\\"\\s]+)\", 1, ProcessCommandLine)\n| extend taskRun  = extract(@\"/tr\\s+\\\"?([^\\\"]+?)\\\"?(?:\\s+/|$)\", 1, ProcessCommandLine)\n| where InitiatingProcessFileName !in~ (\"ccmexec.exe\",\"gpscript.exe\",\"tiworker.exe\")\n| join kind=leftouter (\n    DeviceProcessEvents\n    | where Timestamp > ago(14d)\n    | where FileName in~ (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"wmic.exe\")\n    | where InitiatingProcessFileName =~ \"svchost.exe\" and InitiatingProcessCommandLine has \"Schedule\"\n    | project DeviceId, scriptExec=Timestamp, scriptCmd=ProcessCommandLine\n) on DeviceId\n| where isnull(scriptExec) or scriptExec > Timestamp\n| project Timestamp, DeviceName, AccountName, taskName, taskRun, ProcessCommandLine, scriptCmd",
      "rationale": "Most enterprise scheduled tasks point at signed binaries in Program Files; Higaisa's chain instead writes its stager to AppData/Public/ProgramData and reanimates it via wscript/cscript/mshta or the WMIC XSL trick (T1220) with hidden-window flags. Joining task creation to subsequent svchost->wscript/mshta children initiated by the Schedule service confirms the persistence executed and reduces FPs v"
    }
  ]
}