# intel/.llm_actor_uc_cache/78/784693db970bf0be8a6b0786753da10e40db7f3a.json
{
  "name": "Higaisa",
  "ucs": [
    {
      "title": "[LLM] Higaisa LNK-spearphish chain: cmd.exe + find.exe + mshta self-extracting from .lnk",
      "description": "Hunts Higaisa's signature initial-access tradecraft \u2014 a malicious .lnk attachment whose target line invokes cmd.exe to use find.exe (or findstr.exe) to locate a marker string inside the .lnk itself, then pipes the carved payload to mshta.exe / cscript.exe / powershell.exe for execution. This unusual self-carving trick is documented in the Zscaler, PTSecurity and Malwarebytes 2020 Higaisa write-ups (G0126) and is rarely seen outside this actor.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.003",
        "T1059.005",
        "T1059.007",
        "T1027.013",
        "T1140",
        "T1220"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process) as parent_cmd values(Processes.user) as user values(Processes.process_id) as pid from datamodel=Endpoint.Processes where Processes.parent_process_name=\"cmd.exe\" Processes.process_name IN (\"find.exe\",\"findstr.exe\",\"mshta.exe\",\"cscript.exe\",\"wscript.exe\",\"powershell.exe\") (Processes.parent_process=\"*find*.lnk*\" OR Processes.parent_process=\"*findstr*.lnk*\" OR Processes.parent_process=\"*for /f*find*\" OR Processes.parent_process=\"*mshta*find*\" OR Processes.parent_process=\"*.lnk*mshta*\") by Processes.dest Processes.user Processes.parent_process_name Processes.process_name Processes.process Processes.parent_process | `drop_dm_object_name(Processes)` | where match(parent_cmd,\"(?i)\\.lnk\") AND match(parent_cmd,\"(?i)(find|findstr).*\\.lnk\") | convert ctime(firstTime) ctime(lastTime)",
      "kql": "// Higaisa LNK self-carve chain: cmd.exe spawned by explorer.exe (LNK click), executing find/findstr against a .lnk and piping to a script host\nlet ScriptHosts = dynamic([\"mshta.exe\",\"cscript.exe\",\"wscript.exe\",\"powershell.exe\",\"pwsh.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(30d)\n| where AccountName !endswith \"$\"\n// Stage 1 \u2014 cmd.exe with the find-from-lnk pattern, parent typically explorer.exe (user double-click)\n| where InitiatingProcessFileName =~ \"cmd.exe\"\n| where InitiatingProcessParentFileName in~ (\"explorer.exe\",\"outlook.exe\",\"winrar.exe\",\"7zfm.exe\")\n| where InitiatingProcessCommandLine matches regex @\"(?i)(find|findstr).{1,40}\\.lnk\"\n   or InitiatingProcessCommandLine matches regex @\"(?i)for\\s+/f.{1,80}\\.lnk\"\n// Stage 2 \u2014 child is a script host (mshta/cscript/wscript/powershell)\n| where FileName in~ (ScriptHosts)\n| project Timestamp, DeviceName, AccountName,\n          GrandParent = InitiatingProcessParentFileName,\n          ParentCmd   = InitiatingProcessCommandLine,\n          ChildBin    = FileName,\n          ChildCmd    = ProcessCommandLine,\n          SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// Higaisa LNK self-carve chain \u2014 Sentinel SecurityEvent (EventID 4688) variant\nlet ScriptHosts = dynamic([\"mshta.exe\",\"cscript.exe\",\"wscript.exe\",\"powershell.exe\",\"pwsh.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(30d)\n| where EventID == 4688\n| where AccountType == \"User\"\n| where Account !endswith \"$\"\n| where ParentProcessName endswith \"cmd.exe\"\n| where CommandLine matches regex @\"(?i)(find|findstr).{1,40}\\.lnk\"\n     or CommandLine matches regex @\"(?i)for\\s+/f.{1,80}\\.lnk\"\n// Pivot to the spawned script-host child within ~5s on the same host/session\n| project ParentTime = TimeGenerated, Computer, Account,\n          ParentCmd = CommandLine, SubjectLogonId\n| join kind=inner (\n    SecurityEvent\n    | where TimeGenerated > ago(30d)\n    | where EventID == 4688\n    | where tolower(NewProcessName) has_any (\"\\\\mshta.exe\",\"\\\\cscript.exe\",\"\\\\wscript.exe\",\"\\\\powershell.exe\",\"\\\\pwsh.exe\")\n    | project ChildTime = TimeGenerated, Computer, Account,\n              ChildCmd = CommandLine, ChildBin = NewProcessName, SubjectLogonId\n  ) on Computer, Account, SubjectLogonId\n| where ChildTime between (ParentTime .. ParentTime + 10s)\n| project ParentTime, ChildTime, Computer, Account, ParentCmd, ChildBin, ChildCmd\n| order by ParentTime desc",
      "sigma_yaml": "title: Higaisa LNK Self-Extracting Spearphish (find/findstr + mshta chain)\nid: e8f4c7a1-h1g4-4i5a-b2c3-higaisa-lnk-001\nstatus: experimental\ndescription: Detects the Higaisa (G0126) signature LNK initial-access tradecraft where a .lnk target invokes cmd.exe to find/findstr a marker string within the .lnk itself, then pipes the carved payload to a script host (mshta/cscript/wscript/powershell).\nreferences:\n  - https://attack.mitre.org/groups/G0126/\n  - https://www.zscaler.com/blogs/security-research/return-higaisa-apt\n  - https://www.malwarebytes.com/blog/news/2020/06/higaisa\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.initial_access\n  - attack.execution\n  - attack.t1566.001\n  - attack.t1204.002\n  - attack.t1059.003\n  - attack.g0126\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent_cmd:\n    ParentImage|endswith: '\\cmd.exe'\n    ParentCommandLine|re: '(?i)(find|findstr|for\\s+/f).{1,80}\\.lnk'\n  selection_child_bin:\n    Image|endswith:\n      - '\\mshta.exe'\n      - '\\cscript.exe'\n      - '\\wscript.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n  condition: selection_parent_cmd and selection_child_bin\nfalsepositives:\n  - Niche admin scripts that grep inside .lnk files (very rare)\nlevel: high",
      "rationale": "Higaisa's LNK lures (documented by Zscaler, PTSecurity, Malwarebytes in 2020) are unusual in that the .lnk target line embeds a cmd.exe one-liner that uses find.exe / findstr.exe to locate a marker string INSIDE the .lnk file itself, decode the carved bytes, and feed them to mshta.exe (most often) or another script host. This combination \u2014 cmd.exe spawned from explorer/outlook/archive-tool whose c",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] Higaisa persistence: schtasks creating tasks pointing to %APPDATA%\\Microsoft staged payloads",
      "description": "Hunts Higaisa's documented persistence pattern: after the LNK chain detonates, the actor stages payloads (svchost.exe / SManager.exe / officeupdate.exe lookalikes) under %APPDATA%\\Microsoft\\* and registers a Scheduled Task via schtasks.exe pointing at that path. Tied to G0126 by the consistent %APPDATA%\\Microsoft staging directory and benign-sounding task names observed in the 2020 reports.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1053.005",
        "T1547.001",
        "T1036.004",
        "T1574.001",
        "T1564.003"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" (Processes.process=\"*/create*\" OR Processes.process=\"*-create*\") (Processes.process=\"*AppData\\\\Roaming\\\\Microsoft*\" OR Processes.process=\"*AppData\\\\Local\\\\Microsoft*\" OR Processes.process=\"*\\\\Microsoft\\\\svchost.exe*\" OR Processes.process=\"*SManager*\" OR Processes.process=\"*officeupdate*\" OR Processes.process=\"*\\\\Microsoft\\\\Vault\\\\*\") by Processes.dest Processes.user Processes.parent_process_name Processes.process | `drop_dm_object_name(Processes)` | where NOT match(user,\"(?i)\\\\$$\") AND NOT match(parent,\"(?i)(msiexec|TiWorker|trustedinstaller)\\\\.exe\") | convert ctime(firstTime) ctime(lastTime)",
      "kql": "// Higaisa scheduled-task persistence pointing into %APPDATA%\\Microsoft staging dirs\nDeviceProcessEvents\n| where Timestamp > ago(30d)\n| where AccountName !endswith \"$\"\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has_any (\"/create\",\"-create\")\n// Higaisa-specific staging paths and binary names from Zscaler/Malwarebytes 2020 IoCs\n| where ProcessCommandLine matches regex @\"(?i)appdata\\\\(roaming|local)\\\\microsoft\\\\\"\n     or ProcessCommandLine has_any (\"\\\\Microsoft\\\\svchost.exe\",\"SManager.exe\",\"officeupdate\",\"\\\\Microsoft\\\\Vault\\\\\")\n// Suspicious parent: not a normal installer\n| where InitiatingProcessFileName !in~ (\"msiexec.exe\",\"trustedinstaller.exe\",\"tiworker.exe\")\n| project Timestamp, DeviceName, AccountName,\n          Parent = InitiatingProcessFileName,\n          ParentCmd = InitiatingProcessCommandLine,\n          SchtasksCmd = ProcessCommandLine,\n          SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// Higaisa scheduled-task persistence \u2014 Sentinel 4688 variant\nSecurityEvent\n| where TimeGenerated > ago(30d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| where NewProcessName endswith \"\\\\schtasks.exe\"\n| where CommandLine has_any (\"/create\",\"-create\")\n| where CommandLine matches regex @\"(?i)appdata\\\\(roaming|local)\\\\microsoft\\\\\"\n     or CommandLine has_any (\"\\\\Microsoft\\\\svchost.exe\",\"SManager.exe\",\"officeupdate\",\"\\\\Microsoft\\\\Vault\\\\\")\n| where ParentProcessName !endswith \"msiexec.exe\"\n     and ParentProcessName !endswith \"TrustedInstaller.exe\"\n     and ParentProcessName !endswith \"TiWorker.exe\"\n| project TimeGenerated, Computer, Account, ParentProcessName, NewProcessName, CommandLine\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Higaisa Scheduled Task Persistence to %APPDATA%\\Microsoft\nid: 4d2b8c11-h1g4-4ab2-9f4c-higaisa-sch-002\nstatus: experimental\ndescription: Detects schtasks.exe creating a task whose action targets %APPDATA%\\Microsoft\\ staging directories or filenames (svchost.exe under Microsoft folder, SManager.exe, officeupdate.exe) consistent with Higaisa (G0126) post-exploitation persistence.\nreferences:\n  - https://attack.mitre.org/groups/G0126/\n  - https://www.zscaler.com/blogs/security-research/return-higaisa-apt\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.persistence\n  - attack.t1053.005\n  - attack.t1036.004\n  - attack.g0126\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_schtasks:\n    Image|endswith: '\\schtasks.exe'\n    CommandLine|contains:\n      - '/create'\n      - '-create'\n  selection_higaisa_path:\n    CommandLine|re: '(?i)appdata\\\\(roaming|local)\\\\microsoft\\\\'\n  selection_higaisa_names:\n    CommandLine|contains:\n      - '\\Microsoft\\svchost.exe'\n      - 'SManager.exe'\n      - 'officeupdate'\n      - '\\Microsoft\\Vault\\'\n  filter_legit_parent:\n    ParentImage|endswith:\n      - '\\msiexec.exe'\n      - '\\TrustedInstaller.exe'\n      - '\\TiWorker.exe'\n  condition: selection_schtasks and (selection_higaisa_path or selection_higaisa_names) and not filter_legit_parent\nfalsepositives:\n  - Custom in-house deployment scripts that drop to %APPDATA%\\Microsoft (rare \u2014 usually Program Files / ProgramData)\nlevel: high",
      "rationale": "The combination of (a) schtasks.exe /create, (b) the action path being %APPDATA%\\Microsoft (a folder Microsoft installers themselves never use \u2014 legitimate installers go to %ProgramFiles% or %ProgramData%), and (c) the specific masqueraded filenames svchost.exe, SManager.exe and officeupdate.exe published as IoCs in the Zscaler 2020 Higaisa report, makes this a tight fit for G0126's post-LNK persi",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    }
  ]
}