# intel/.llm_actor_uc_cache/81/8146f2f2116f02e4ccd81260717e168d3df23ce4.json
{
  "name": "Naikon",
  "ucs": [
    {
      "title": "[LLM] Naikon (PLA Unit 78020) Aria-body / RainyDay DLL search-order hijack via signed-binary side-load",
      "description": "Hunts Naikon's signature delivery technique: dropping a legitimately-signed Microsoft/AV/VPN executable into AppData, Temp, Public or ProgramData and using DLL search-order hijacking so the executable side-loads a malicious DLL (Aria-body, RainyDay, Nebulae) from the same folder. CrowdStrike tracks Naikon as part of the Lotus Panda/Override Panda activity cluster; Kaspersky and Bitdefender both documented this exact loader pattern across SE-Asian government and military intrusions (2015-2021).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.001",
        "T1036.005",
        "T1204.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstSeen max(_time) as lastSeen values(Processes.process) as cmdline values(Processes.process_hash) as parent_hash from datamodel=Endpoint.Processes where Processes.process_path IN (\"*\\\\AppData\\\\*\",\"*\\\\Temp\\\\*\",\"*\\\\Users\\\\Public\\\\*\",\"*\\\\ProgramData\\\\*\") AND (Processes.process_company=\"Microsoft Corporation\" OR Processes.process_company=\"Kaspersky*\" OR Processes.process_company=\"ESET*\" OR Processes.process_company=\"Symantec*\" OR Processes.process_company=\"VMware*\" OR Processes.process_company=\"Adobe*\") by host Processes.user Processes.process_name Processes.process_path Processes.process_hash\n| `drop_dm_object_name(Processes)`\n| join type=inner host process_path [\n    | tstats summariesonly=true count from datamodel=Endpoint.Filesystem where Filesystem.file_name=\"*.dll\" AND NOT Filesystem.file_path IN (\"C:\\\\Windows\\\\*\",\"C:\\\\Program Files*\\\\*\") by host Filesystem.file_path Filesystem.file_name Filesystem.file_hash\n    | `drop_dm_object_name(Filesystem)`\n    | rex field=file_path \"(?<dll_dir>.*)\\\\\\\\[^\\\\\\\\]+$\"\n    | rename dll_dir as process_path file_hash as dll_hash file_name as dll_name\n]\n| where dll_hash != process_hash\n| table firstSeen host user process_name process_path process_hash dll_name dll_hash cmdline",
      "kql": "// Naikon Aria-body / RainyDay side-load: signed exe in user-writeable path\n// loads a non-system DLL from the SAME folder.\nlet LookbackDays = 14d;\nlet _user_writeable = dynamic([@\"\\AppData\\\", @\"\\Temp\\\", @\"\\Users\\Public\\\", @\"\\ProgramData\\\"]);\nlet _signed_companies = dynamic([\"Microsoft Corporation\",\"Kaspersky Lab\",\"ESET\",\"ESET, spol. s r.o.\",\"Symantec Corporation\",\"VMware, Inc.\",\"Adobe Inc.\",\"Adobe Systems Incorporated\",\"Avast Software\",\"Outlook\"]);\nDeviceImageLoadEvents\n| where Timestamp > ago(LookbackDays)\n| where InitiatingProcessFolderPath has_any (_user_writeable)              // dropped to a writeable path \u2014 Naikon TTP\n| where InitiatingProcessVersionInfoCompanyName in~ (_signed_companies)    // legit signed loader binary\n| where FileName endswith \".dll\"\n| where not(FolderPath startswith @\"C:\\Windows\\\")                         // DLL not from system\n| where not(FolderPath startswith @\"C:\\Program Files\")\n| where SHA256 != InitiatingProcessSHA256                                  // distinct file (not the exe itself)\n| extend ExeDir = tolower(InitiatingProcessFolderPath)\n| extend DllDir = tolower(FolderPath)\n| where DllDir == ExeDir                                                   // DLL planted next to the exe\n| project Timestamp, DeviceName,\n          LoaderExe = InitiatingProcessFileName,\n          LoaderPath = InitiatingProcessFolderPath,\n          LoaderSigner = InitiatingProcessVersionInfoCompanyName,\n          LoaderSHA256 = InitiatingProcessSHA256,\n          SideLoadedDll = FileName,\n          SideLoadedSHA256 = SHA256,\n          InitiatingProcessAccountName, InitiatingProcessCommandLine\n| order by Timestamp desc",
      "sentinel_kql": "// Requires Sysmon EID 7 (Image Loaded) ingested via WindowsEvent.\nlet LookbackDays = 14d;\nlet _writeable_rx = @\"(?i)\\\\(AppData|Temp|Public|ProgramData)\\\\\";\nlet _signed = dynamic([\"Microsoft Corporation\",\"Kaspersky Lab\",\"ESET\",\"Symantec Corporation\",\"VMware, Inc.\",\"Adobe Inc.\",\"Avast Software\"]);\nWindowsEvent\n| where TimeGenerated > ago(LookbackDays)\n| where Provider == \"Microsoft-Windows-Sysmon\" and EventID == 7\n| extend Image = tostring(EventData.Image),\n         ImageLoaded = tostring(EventData.ImageLoaded),\n         LoaderCompany = tostring(EventData.Company),\n         LoaderSigned = tostring(EventData.Signed),\n         LoadedHash = tostring(EventData.Hashes),\n         User = tostring(EventData.User)\n| where Image matches regex _writeable_rx\n| where LoaderCompany in~ (_signed)\n| where ImageLoaded endswith \".dll\"\n| where not(ImageLoaded startswith @\"C:\\Windows\\\") and not(ImageLoaded startswith @\"C:\\Program Files\")\n| extend ExeDir = tolower(strcat_array(array_slice(split(Image, @\"\\\"), 0, -2), @\"\\\")),\n         DllDir = tolower(strcat_array(array_slice(split(ImageLoaded, @\"\\\"), 0, -2), @\"\\\"))\n| where ExeDir == DllDir\n| project TimeGenerated, Computer, User, LoaderExe = Image, LoaderCompany, LoaderSigned,\n          SideLoadedDll = ImageLoaded, LoadedHash\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Naikon Aria-body / RainyDay Signed-Binary DLL Side-Load\nid: a3d7f8b2-9c11-4e8f-9b6a-7c4d8e2a1f30\nstatus: experimental\ndescription: Detects legitimately-signed loader binaries (Microsoft, Kaspersky, ESET, Symantec, VMware, Adobe) executing from user-writeable paths and loading a non-system DLL from the same directory. This is the canonical Naikon (Lotus Panda, PLA Unit 78020) Aria-body / RainyDay / Nebulae loader pattern.\nreferences:\n  - https://attack.mitre.org/groups/G0019/\n  - https://securelist.com/the-naikon-apt/\n  - https://www.bitdefender.com/files/News/CaseStudies/study/388/Bitdefender-PR-Whitepaper-NAIKON-creat5397-en-EN.pdf\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.t1574.001\n  - attack.t1036.005\n  - attack.persistence\n  - attack.defense_evasion\nlogsource:\n  product: windows\n  category: image_load\ndetection:\n  selection_path:\n    Image|contains:\n      - '\\AppData\\'\n      - '\\Temp\\'\n      - '\\Users\\Public\\'\n      - '\\ProgramData\\'\n  selection_signed:\n    Company:\n      - 'Microsoft Corporation'\n      - 'Kaspersky Lab'\n      - 'ESET'\n      - 'Symantec Corporation'\n      - 'VMware, Inc.'\n      - 'Adobe Inc.'\n      - 'Avast Software'\n  selection_dll:\n    ImageLoaded|endswith: '.dll'\n  filter_system:\n    ImageLoaded|startswith:\n      - 'C:\\Windows\\'\n      - 'C:\\Program Files'\n  condition: all of selection_* and not filter_system\nfields:\n  - Image\n  - ImageLoaded\n  - Company\n  - User\nfalsepositives:\n  - Portable installers that legitimately ship signed binaries with bundled DLLs\n  - Software updaters staging files in ProgramData\nlevel: high",
      "rationale": "Naikon's defining tradecraft across every public report (Kaspersky 2015 'CameraShy', Bitdefender 2020 Aria-body, ESET 2022 Nebulae) is signed-binary DLL side-loading from user-writeable paths. Generic side-load detections fire on installers; pinning the parent's signer to a small list of impersonated vendors AND requiring the DLL to sit next to the exe is what makes this Naikon-specific rather tha",
      "_field_issues": [
        {
          "kind": "wrong_table",
          "field": "InitiatingProcessVersionInfoCompanyName",
          "tables_in_scope": [
            "DeviceImageLoadEvents"
          ],
          "exists_on": [
            "DeviceProcessEvents",
            "DeviceNetworkEvents",
            "DeviceFileEvents"
          ],
          "message": "`InitiatingProcessVersionInfoCompanyName` is not a column of DeviceImageLoadEvents; it does exist on DeviceProcessEvents, DeviceNetworkEvents, DeviceFileEvents."
        }
      ]
    },
    {
      "title": "[LLM] Naikon WMIC reconnaissance burst followed by scheduled-task or WMI persistence within 30 min",
      "description": "Naikon operators repeatedly run a stereotyped WMIC discovery sequence (process list, antivirus product enumeration, OS/computersystem, NIC config, domain accounts) immediately before installing persistence via schtasks.exe or 'wmic process call create'. Catching the recon-to-persistence chain on a single host within 30 minutes \u2014 by the same actor account \u2014 is far higher fidelity than alerting on either side alone, and matches the playbook documented in the Bitdefender/ESET intrusion timelines.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1047",
        "T1518.001",
        "T1053.005",
        "T1016",
        "T1018",
        "T1078.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as recon_time values(Processes.process) as recon_cmd from datamodel=Endpoint.Processes where Processes.process_name=\"wmic.exe\" AND (Processes.process=\"*process list*\" OR Processes.process=\"*process get*\" OR Processes.process=\"*/node:*\" OR Processes.process=\"*antivirus*\" OR Processes.process=\"*computersystem*\" OR Processes.process=\"*useraccount*\" OR Processes.process=\"*ntdomain*\" OR Processes.process=\"*nicconfig*\") by host Processes.user\n| `drop_dm_object_name(Processes)`\n| rename user as recon_user\n| join type=inner host [\n    | tstats summariesonly=true count min(_time) as persist_time values(Processes.process) as persist_cmd from datamodel=Endpoint.Processes where (Processes.process_name=\"schtasks.exe\" AND Processes.process=\"*/create*\") OR (Processes.process_name=\"wmic.exe\" AND Processes.process=\"*process*\" AND Processes.process=\"*call*\" AND Processes.process=\"*create*\") by host Processes.user\n    | `drop_dm_object_name(Processes)`\n    | rename user as persist_user\n]\n| eval delta_sec=persist_time-recon_time\n| where delta_sec>=0 AND delta_sec<=1800 AND recon_user=persist_user AND NOT match(recon_user,\"\\\\$$\")\n| table recon_time persist_time delta_sec host recon_user recon_cmd persist_cmd\n| sort - recon_time",
      "kql": "// Naikon recon \u2192 persistence chain: WMIC discovery cluster followed by\n// schtasks /create or 'wmic process call create' on the same host\n// by the same account within 30 minutes.\nlet LookbackHours = 24h;\nlet ChainWindow = 30m;\nlet _wmic_recon_tokens = dynamic([\n    \"process list\",\"process get\",\"/node:\",\"antivirus\",\"computersystem\",\n    \"useraccount\",\"ntdomain\",\"nicconfig\",\"product get\",\"service get\",\"os get\"\n]);\nlet WmicRecon =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackHours)\n    | where FileName =~ \"wmic.exe\"\n    | where ProcessCommandLine has_any (_wmic_recon_tokens)\n    | where AccountName !endswith \"$\"                                 // skip machine accounts\n    | summarize ReconStart = min(Timestamp),\n                ReconCmds = make_set(ProcessCommandLine, 25),\n                ReconCount = count()\n        by DeviceId, DeviceName, AccountName\n    | where ReconCount >= 2;                                          // require a CLUSTER, not a one-off\nDeviceProcessEvents\n| where Timestamp > ago(LookbackHours)\n| where AccountName !endswith \"$\"\n| where (FileName =~ \"schtasks.exe\" and ProcessCommandLine has \"/create\")\n     or (FileName =~ \"wmic.exe\" and ProcessCommandLine has_all (\"process\",\"call\",\"create\"))\n| project PersistTime = Timestamp, DeviceId, DeviceName, AccountName,\n          PersistBin = FileName,\n          PersistCmd = ProcessCommandLine,\n          PersistParent = InitiatingProcessFileName,\n          PersistParentCmd = InitiatingProcessCommandLine\n| join kind=inner WmicRecon on DeviceId, AccountName\n| where PersistTime between (ReconStart .. ReconStart + ChainWindow)\n| extend DeltaSec = datetime_diff('second', PersistTime, ReconStart)\n| project ReconStart, PersistTime, DeltaSec, DeviceName, AccountName,\n          ReconCount, ReconCmds, PersistBin, PersistCmd, PersistParent, PersistParentCmd\n| order by ReconStart desc",
      "sentinel_kql": "// Sentinel via 4688 (requires 'Include command line in process creation events' GPO).\nlet LookbackHours = 24h;\nlet ChainWindow = 30m;\nlet _wmic_recon_tokens = dynamic([\n    \"process list\",\"process get\",\"/node:\",\"antivirus\",\"computersystem\",\n    \"useraccount\",\"ntdomain\",\"nicconfig\",\"product get\",\"service get\",\"os get\"\n]);\nlet WmicRecon =\n    SecurityEvent\n    | where TimeGenerated > ago(LookbackHours)\n    | where EventID == 4688\n    | where NewProcessName endswith @\"\\wmic.exe\"\n    | where CommandLine has_any (_wmic_recon_tokens)\n    | where Account !endswith \"$\"\n    | summarize ReconStart = min(TimeGenerated),\n                ReconCmds = make_set(CommandLine, 25),\n                ReconCount = count()\n        by Computer, Account\n    | where ReconCount >= 2;\nSecurityEvent\n| where TimeGenerated > ago(LookbackHours)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| where (NewProcessName endswith @\"\\schtasks.exe\" and CommandLine has \"/create\")\n     or (NewProcessName endswith @\"\\wmic.exe\" and CommandLine has_all (\"process\",\"call\",\"create\"))\n| project PersistTime = TimeGenerated, Computer, Account,\n          PersistBin = NewProcessName, PersistCmd = CommandLine,\n          PersistParent = ParentProcessName\n| join kind=inner WmicRecon on Computer, Account\n| where PersistTime between (ReconStart .. ReconStart + ChainWindow)\n| extend DeltaSec = datetime_diff('second', PersistTime, ReconStart)\n| project ReconStart, PersistTime, DeltaSec, Computer, Account,\n          ReconCount, ReconCmds, PersistBin, PersistCmd, PersistParent\n| order by ReconStart desc",
      "sigma_yaml": "",
      "rationale": "Each individual leg (WMIC recon, schtasks /create, 'wmic process call create') is too noisy for an alert in a real estate \u2014 sysadmins, SCCM, and patching tools generate hundreds per day. Naikon is one of the very few actors who routinely chains a *cluster* of WMIC discovery commands (especially the AV/computersystem/ntdomain combo) directly into persistence via schtasks or remote WMI process creat"
    }
  ]
}