# intel/.llm_actor_uc_cache/81/8146f2f2116f02e4ccd81260717e168d3df23ce4.json
{
  "name": "Naikon",
  "ucs": [
    {
      "title": "[LLM] Naikon (Override Panda) DLL search-order hijack via signed AV/Office binaries loading Aria-body / RainyDay / Nebulae",
      "description": "Hunts Naikon's signature loader pattern: legitimately-signed binaries (Outlook, Kaspersky avp.exe, Avast wsc_proxy, Sandboxie, Symantec, Mcafee mcods.exe) sideloading attacker-supplied DLLs (e.g. aspnet_perf.dll, mscoree.dll, oci.dll, k7sysmon.exe.dll, fltlib.dll) staged in user-writable paths to inject the Aria-body, RainyDay or Nebulae backdoors used against ASEAN government and military targets.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.001",
        "T1036.005",
        "T1218"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process) as parent_cmdline from datamodel=Endpoint.Processes where Processes.process_name IN (\"outlook.exe\",\"avp.exe\",\"wsc_proxy.exe\",\"mcods.exe\",\"SandboxieRpcSs.exe\",\"sandman.exe\",\"k7sysmon.exe\",\"mscorsvw.exe\",\"AdobeUpdate.exe\") by host Processes.process_name Processes.process_path Processes.user _time span=1h \n| join type=inner host [| tstats summariesonly=t values(Filesystem.file_path) as dll_path values(Filesystem.file_name) as dll_name from datamodel=Endpoint.Filesystem where Filesystem.file_name IN (\"aspnet_perf.dll\",\"mscoree.dll\",\"oci.dll\",\"fltlib.dll\",\"version.dll\",\"DismCore.dll\",\"k7sysmon.exe.dll\",\"Sandboxie.dll\",\"Outllib.dll\",\"libcurl.dll\",\"msvcr100.dll\") by host Filesystem.file_path] \n| where like(Processes.process_path, \"%\\\\Users\\\\%\") OR like(Processes.process_path, \"%\\\\ProgramData\\\\%\") OR like(Processes.process_path, \"%\\\\AppData\\\\%\") OR like(Processes.process_path, \"%\\\\Public\\\\%\") OR like(dll_path, \"%\\\\Users\\\\%\") OR like(dll_path, \"%\\\\ProgramData\\\\%\") \n| eval naikon_indicator=if((Processes.process_name=\"outlook.exe\" AND match(dll_name,\"Outllib|mscoree\")) OR (Processes.process_name=\"avp.exe\" AND match(dll_name,\"oci|fltlib\")) OR (Processes.process_name=\"k7sysmon.exe\" AND dll_name=\"k7sysmon.exe.dll\"),\"high\",\"medium\") \n| table firstTime lastTime host Processes.user Processes.process_name Processes.process_path dll_name dll_path naikon_indicator",
      "kql": "let naikon_hosts = dynamic([\"outlook.exe\",\"avp.exe\",\"wsc_proxy.exe\",\"mcods.exe\",\"SandboxieRpcSs.exe\",\"sandman.exe\",\"k7sysmon.exe\",\"mscorsvw.exe\",\"AdobeUpdate.exe\",\"smtools.exe\"]);\nlet naikon_dlls  = dynamic([\"aspnet_perf.dll\",\"mscoree.dll\",\"oci.dll\",\"fltlib.dll\",\"version.dll\",\"DismCore.dll\",\"k7sysmon.exe.dll\",\"Sandboxie.dll\",\"Outllib.dll\",\"libcurl.dll\",\"msvcr100.dll\"]);\nlet writable_paths = dynamic([@\"\\users\\\",@\"\\programdata\\\",@\"\\appdata\\\",@\"\\public\\\",@\"\\temp\\\",@\"\\perflogs\\\"]);\nDeviceImageLoadEvents\n| where InitiatingProcessFileName in~ (naikon_hosts)\n| where FileName in~ (naikon_dlls)\n| extend FolderLower = tolower(FolderPath), HostLower = tolower(InitiatingProcessFolderPath)\n| where FolderLower has_any (writable_paths) or HostLower has_any (writable_paths)\n| extend signed_host = isnotempty(InitiatingProcessSignerType) and InitiatingProcessSignatureStatus == \"Valid\"\n| where signed_host == true\n| extend NaikonHighFidelity = case(\n    InitiatingProcessFileName =~ \"outlook.exe\" and FileName has_any (\"Outllib\",\"mscoree\"), \"AriaBody/Nebulae loader\",\n    InitiatingProcessFileName =~ \"avp.exe\" and FileName has_any (\"oci\",\"fltlib\"), \"RainyDay loader\",\n    InitiatingProcessFileName =~ \"k7sysmon.exe\" and FileName =~ \"k7sysmon.exe.dll\", \"RainyDay K7 sideload\",\n    InitiatingProcessFileName =~ \"mcods.exe\", \"Aria-body McAfee sideload\",\n    \"Possible Naikon sideload\")\n| project Timestamp, DeviceName, AccountName=InitiatingProcessAccountName, InitiatingProcessFileName, InitiatingProcessFolderPath, FileName, FolderPath, SHA256, NaikonHighFidelity\n| join kind=leftouter (DeviceProcessEvents | where ProcessCommandLine has_any (\"\\\\Users\\\\Public\",\"\\\\ProgramData\\\\\",\"rundll32\",\"regsvr32\") | project DeviceName, ProcessCommandLine, Timestamp) on DeviceName\n| sort by Timestamp desc",
      "rationale": "DLL search-order hijacking against AV/Office binaries to load Aria-body (Naikon's flagship RAT documented by Check Point in 2020), RainyDay and Nebulae is the defining tradecraft of Naikon / G0019 / Override Panda. The specific pairings (avp.exe+oci.dll, k7sysmon.exe+k7sysmon.exe.dll, outlook.exe+Outllib.dll, mcods.exe sideload) are unique fingerprints documented in Check Point, Bitdefender and Ka"
    },
    {
      "title": "[LLM] Naikon RoyalRoad RTF lure \u2192 schtasks persistence \u2192 internal recon chain (G0019 spearphishing playbook)",
      "description": "Catches Naikon's end-to-end intrusion chain against Southeast-Asian diplomatic and military targets: weaponised RTF (RoyalRoad/8.t builder shared with TA428/Goblin Panda) delivered by spearphish, dropping a payload that registers a Scheduled Task and immediately performs SMB/AD reconnaissance (net view, nbtstat, T1018/T1046) \u2014 the same sequence observed in Naikon's 'Aria-body' campaigns against ASEAN ministries.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "medium",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1053.005",
        "T1018",
        "T1046"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count from datamodel=Email where Email.file_name IN (\"*.rtf\",\"*.doc\",\"*.docx\") by Email.recipient Email.src_user Email.subject Email.message_id _time span=1h \n| join type=inner Email.recipient [| tstats summariesonly=t values(Processes.process) as proc_cmdline values(Processes.process_name) as proc_name from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"winword.exe\",\"eqnedt32.exe\") (Processes.process_name IN (\"schtasks.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\") OR Processes.process IN (\"*\\\\8.t*\",\"*Intel*.exe*\",\"*RasTls*\")) by host Processes.user _rename Email.recipient=Processes.user] \n| append [| tstats summariesonly=t count from datamodel=Endpoint.Processes where Processes.process_name IN (\"net.exe\",\"net1.exe\",\"nbtstat.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"netstat.exe\") (Processes.process IN (\"*view /domain*\",\"*group \\\"domain admins\\\"*\",\"*-A *\",\"* /all*\")) by host Processes.user Processes.parent_process_name _time span=10m] \n| stats values(Email.subject) as subject values(proc_name) as proc_chain values(Processes.process) as recon_cmds count by host \n| where count >= 3 \n| eval naikon_score = if(match(mvjoin(proc_chain,\",\"),\"eqnedt32|schtasks\") AND mvcount(recon_cmds)>=2,\"high\",\"medium\")",
      "kql": "let lookback = 24h;\nlet rtf_recipients =\n    EmailAttachmentInfo\n    | where Timestamp > ago(lookback)\n    | where FileName has_any (\".rtf\",\".doc\",\".docx\")\n    | join kind=inner (EmailEvents | where Timestamp > ago(lookback) and DeliveryAction == \"Delivered\") on NetworkMessageId\n    | project Timestamp, RecipientEmailAddress, SenderFromAddress, Subject, FileName, SHA256, NetworkMessageId;\nlet rr_execution =\n    DeviceProcessEvents\n    | where Timestamp > ago(lookback)\n    | where InitiatingProcessFileName in~ (\"winword.exe\",\"eqnedt32.exe\",\"excel.exe\")\n    | where FileName in~ (\"schtasks.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"wmic.exe\")\n        or ProcessCommandLine has_any (\"8.t\",\"intel.wll\",\"RasTls\",\"\\\\Intel\\\\\",\"\\\\Public\\\\\",\"aria-body\",\"AppData\\\\Roaming\\\\Intel\")\n    | project ExecTime=Timestamp, DeviceName, AccountUpn=InitiatingProcessAccountUpn, InitiatingProcessFileName, FileName, ProcessCommandLine;\nlet recon =\n    DeviceProcessEvents\n    | where Timestamp > ago(lookback)\n    | where FileName in~ (\"net.exe\",\"net1.exe\",\"nbtstat.exe\",\"systeminfo.exe\",\"tasklist.exe\",\"netstat.exe\",\"ipconfig.exe\",\"whoami.exe\",\"ping.exe\")\n    | where ProcessCommandLine has_any (\"view /domain\",\"group \\\"domain admins\\\"\",\"/all\",\"-A \",\"-an\",\"/domain\")\n    | summarize ReconCount=dcount(FileName), ReconCmds=make_set(ProcessCommandLine,10), FirstRecon=min(Timestamp), LastRecon=max(Timestamp) by DeviceName, AccountUpn=InitiatingProcessAccountUpn;\nrr_execution\n| join kind=inner rtf_recipients on $left.AccountUpn == $right.RecipientEmailAddress\n| where ExecTime between (Timestamp .. Timestamp + 2h)\n| join kind=inner recon on DeviceName, AccountUpn\n| where FirstRecon between (ExecTime .. ExecTime + 1h) and ReconCount >= 2\n| extend NaikonChain=\"Phish RTF -> Office child -> schtasks/persist -> AD recon\"\n| project Timestamp, SenderFromAddress, Subject, FileName, SHA256, DeviceName, AccountUpn, InitiatingProcessFileName, ProcessCommandLine, ReconCmds, NaikonChain\n| sort by Timestamp desc",
      "rationale": "Naikon's documented playbook (CameraShy report, Check Point 'Naikon's New Backdoor', Bitdefender 2020) consistently chains RoyalRoad-style RTF lures (CVE-2017-11882 via eqnedt32.exe, or .wll Word add-in T1137.006) into immediate scheduled-task persistence and aggressive 'net view /domain' / nbtstat AD enumeration to map ASEAN ministry networks before deploying Aria-body. Joining the email recipien"
    }
  ]
}