# intel/.llm_actor_uc_cache/04/041edba9a3c8b9487c3f6ccb62a54560c13a6478.json
{
  "name": "Patchwork",
  "ucs": [
    {
      "title": "[LLM] Patchwork (Dropping Elephant) Equation Editor exploit chain — EQNEDT32.exe spawns scripting / LOLBin child",
      "description": "Hunts Patchwork's signature delivery chain: a weaponised RTF/DOC lure exploiting CVE-2017-11882 / CVE-2018-0802 in Microsoft Equation Editor (EQNEDT32.exe), which then spawns a scripting host (mshta/wscript/powershell/cmd) or LOLBin (mshta/regsvr32/certutil) to stage BADNEWS, NDiskMonitor or QuasarRAT. This is the most consistently observed initial-access pattern in Patchwork campaigns (Operation Hangover, MONSOON, BADNEWS) since 2017 and is essentially never seen as benign behaviour.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "exploit",
      "conf": "high",
      "techs": [
        "T1203",
        "T1566.001",
        "T1204.002",
        "T1059.001",
        "T1059.005",
        "T1218"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.process_hash) as hashes values(Processes.user) as user values(Processes.parent_process) as parent_cmd from datamodel=Endpoint.Processes where Processes.parent_process_name=\"EQNEDT32.EXE\" Processes.process_name IN (\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"hh.exe\",\"schtasks.exe\") by host Processes.parent_process_name Processes.process_name Processes.user _time span=1m | `drop_dm_object_name(Processes)` | rename firstTime as first_seen lastTime as last_seen | convert ctime(first_seen) ctime(last_seen) | sort - last_seen",
      "kql": "// Patchwork (G0040) — Equation Editor exploit chain (CVE-2017-11882 / CVE-2018-0802)\nlet _suspicious_children = dynamic([\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"hh.exe\",\"schtasks.exe\",\"forfiles.exe\"]);\nDeviceProcessEvents\n| where Timestamp > ago(30d)\n| where InitiatingProcessFileName =~ \"eqnedt32.exe\"   // Equation Editor — Patchwork's preferred CVE-2017-11882 trigger\n| where AccountName !endswith \"$\"\n| where FileName in~ (_suspicious_children)\n   or ProcessCommandLine has_any (\"http://\",\"https://\",\"-EncodedCommand\",\"DownloadString\",\"AppData\\\\Local\\\\Temp\",\"\\\\ProgramData\\\\\",\"-w hidden\",\"/c start\")\n| project Timestamp, DeviceName, AccountName,\n          ParentImage = InitiatingProcessFolderPath,\n          ParentCmd   = InitiatingProcessCommandLine,\n          ParentSHA256= InitiatingProcessSHA256,\n          ChildImage  = FolderPath,\n          ChildCmd    = ProcessCommandLine,\n          ChildSHA256 = SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// Patchwork (G0040) — Equation Editor exploit chain on Sentinel SecurityEvent (4688)\nlet _suspicious_children = dynamic([\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"hh.exe\",\"schtasks.exe\",\"forfiles.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(30d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| where ParentProcessName endswith \"\\\\EQNEDT32.EXE\"          // Patchwork CVE-2017-11882 trigger\n| extend ChildBin = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n| where ChildBin in (_suspicious_children)\n   or CommandLine has_any (\"http://\",\"https://\",\"-EncodedCommand\",\"DownloadString\",\"AppData\\\\Local\\\\Temp\",\"\\\\ProgramData\\\\\",\"-w hidden\",\"/c start\")\n| project TimeGenerated, Computer, Account,\n          ParentProcessName, NewProcessName, CommandLine,\n          TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Patchwork APT - Equation Editor Spawning Scripting/LOLBin Child\nid: 8b1a4e1d-1f4a-4f2c-9c0a-7e9bf2c2a401\nstatus: experimental\ndescription: Detects CVE-2017-11882/CVE-2018-0802 exploitation chain associated with Patchwork (Dropping Elephant, MONSOON, G0040). Equation Editor (EQNEDT32.EXE) should never legitimately launch a scripting host or LOLBin.\nreferences:\n  - https://attack.mitre.org/groups/G0040/\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.initial_access\n  - attack.execution\n  - attack.t1203\n  - attack.t1566.001\n  - attack.g0040\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent:\n    ParentImage|endswith: '\\EQNEDT32.EXE'\n  selection_child:\n    Image|endswith:\n      - '\\mshta.exe'\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\cmd.exe'\n      - '\\rundll32.exe'\n      - '\\regsvr32.exe'\n      - '\\certutil.exe'\n      - '\\bitsadmin.exe'\n      - '\\hh.exe'\n      - '\\schtasks.exe'\n  condition: selection_parent and selection_child\nfalsepositives:\n  - None expected; Equation Editor has been deprecated and should not spawn child processes.\nlevel: high",
      "rationale": "Patchwork has been documented by Forcepoint (MONSOON), Cymmetria (Operation BugDrop / Hangover), Trend Micro and Volexity using weaponised RTFs/DOCs that trigger CVE-2017-11882 in Equation Editor as the first-stage code-execution primitive across nearly every campaign from 2017 onward. EQNEDT32.exe spawning ANY child process is anomalous (the editor doesn't legitimately do this), and pairing the p"
    },
    {
      "title": "[LLM] Patchwork BADNEWS dead-drop resolver beacon — non-browser HTTP fetch from blogspot/forum/pastebin",
      "description": "Hunts the BADNEWS / NDiskMonitor C2 retrieval pattern that has defined Patchwork tradecraft for nearly a decade: the implant calls out to a hard-coded set of legitimate web services (Blogspot, WordPress, Forum posts, GitHub raw, Pastebin, feedburner) to fetch an encoded blob containing the real C2 address (T1102.001 dead-drop resolver). The detection finds non-browser, non-Office processes (LOLBins, rundll32, schtasks-spawned binaries, AppData-resident binaries) reaching these specific resolver hosts.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "medium",
      "techs": [
        "T1102.001",
        "T1071.001",
        "T1547.001",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count values(Web.url) as urls values(Web.dest) as dest values(Web.user) as user min(_time) as firstTime max(_time) as lastTime from datamodel=Web where (Web.url=\"*blogspot.*\" OR Web.url=\"*wordpress.com*\" OR Web.url=\"*pastebin.com*\" OR Web.url=\"*paste.ee*\" OR Web.url=\"*feedburner.com*\" OR Web.url=\"*raw.githubusercontent.com*\" OR Web.url=\"*forum*\") Web.app NOT IN (\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"iexplore.exe\",\"brave.exe\",\"opera.exe\",\"safari.exe\",\"outlook.exe\",\"teams.exe\",\"slack.exe\") by host Web.app Web.user Web.url _time span=5m | `drop_dm_object_name(Web)` | rename firstTime as first_seen lastTime as last_seen | convert ctime(first_seen) ctime(last_seen) | where count >= 2 | sort - last_seen",
      "kql": "// Patchwork (G0040) — BADNEWS/NDiskMonitor dead-drop resolver beacon (T1102.001)\nlet _resolver_hosts = dynamic([\"blogspot.\",\"blogger.com\",\"wordpress.com\",\"pastebin.com\",\"paste.ee\",\"feedburner.com\",\"raw.githubusercontent.com\",\"gist.githubusercontent.com\",\"sites.google.com\"]);\nlet _browser_procs = dynamic([\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"iexplore.exe\",\"brave.exe\",\"opera.exe\",\"outlook.exe\",\"teams.exe\",\"slack.exe\",\"acrord32.exe\"]);\nlet _suspicious_user_paths = dynamic([@\"\\AppData\\Roaming\\\",@\"\\AppData\\Local\\Temp\\\",@\"\\ProgramData\\\",@\"\\Public\\\",@\"\\Users\\Public\\\"]);\nDeviceNetworkEvents\n| where Timestamp > ago(14d)\n| where RemoteIPType == \"Public\"\n| where isnotempty(RemoteUrl)\n| where InitiatingProcessAccountName !endswith \"$\"\n| where InitiatingProcessFileName !in~ (_browser_procs)\n| where RemoteUrl has_any (_resolver_hosts)\n| where InitiatingProcessFileName in~ (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"curl.exe\",\"wget.exe\")\n   or InitiatingProcessFolderPath has_any (_suspicious_user_paths)\n| summarize HitCount = count(),\n            FirstSeen = min(Timestamp), LastSeen = max(Timestamp),\n            UrlSamples = make_set(RemoteUrl, 10),\n            CmdSamples = make_set(InitiatingProcessCommandLine, 5)\n            by DeviceName, InitiatingProcessAccountName, InitiatingProcessFileName, InitiatingProcessFolderPath, InitiatingProcessSHA256\n| order by LastSeen desc",
      "sentinel_kql": "// Patchwork (G0040) — dead-drop resolver beacon, ASIM cross-vendor\nlet _resolver_hosts = dynamic([\"blogspot.\",\"blogger.com\",\"wordpress.com\",\"pastebin.com\",\"paste.ee\",\"feedburner.com\",\"raw.githubusercontent.com\",\"gist.githubusercontent.com\",\"sites.google.com\"]);\nlet _browser_procs = dynamic([\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"iexplore.exe\",\"brave.exe\",\"opera.exe\",\"outlook.exe\",\"teams.exe\",\"slack.exe\",\"acrord32.exe\"]);\nlet _suspicious_user_paths = dynamic([@\"\\AppData\\Roaming\\\",@\"\\AppData\\Local\\Temp\\\",@\"\\ProgramData\\\",@\"\\Public\\\"]);\nImWebSession\n| where TimeGenerated > ago(14d)\n| where isnotempty(Url)\n| where Url has_any (_resolver_hosts)\n| extend ProcLower = tolower(ActingProcessName)\n| where isnotempty(ProcLower)\n| where ProcLower !in~ (_browser_procs)\n| where SrcUsername !endswith \"$\"\n| where ProcLower in~ (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"mshta.exe\",\"wscript.exe\",\"cscript.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"certutil.exe\",\"bitsadmin.exe\",\"curl.exe\",\"wget.exe\")\n   or ActingProcessCommandLine has_any (_suspicious_user_paths)\n| summarize HitCount = count(),\n            FirstSeen = min(TimeGenerated), LastSeen = max(TimeGenerated),\n            UrlSamples = make_set(Url, 10),\n            CmdSamples = make_set(ActingProcessCommandLine, 5)\n            by DvcHostname, SrcUsername, ActingProcessName\n| order by LastSeen desc",
      "sigma_yaml": "",
      "rationale": "Patchwork's BADNEWS family (documented by Cymmetria 'Unfinished Business', Forcepoint MONSOON, Volexity, Trend Micro) is famous for hiding the real C2 address inside a comment / encoded blob on Blogspot, WordPress, GitHub gists and forum posts — a deliberate dead-drop resolver pattern (T1102.001). Patchwork's tooling characteristically drops the implant under %APPDATA% or %ProgramData% and beacons",
      "_sentinel_field_issues": [
        {
          "kind": "wrong_table",
          "field": "ActingProcessName",
          "tables_in_scope": [
            "ImWebSession"
          ],
          "exists_on": [
            "ImProcessCreate",
            "ImFileEvent",
            "ImRegistryEvent"
          ],
          "message": "`ActingProcessName` is not a column of ImWebSession; it does exist on ImProcessCreate, ImFileEvent, ImRegistryEvent."
        },
        {
          "kind": "wrong_table",
          "field": "ActorUsername",
          "tables_in_scope": [
            "ImWebSession"
          ],
          "exists_on": [
            "ImProcessCreate",
            "ImAuthentication",
            "ImFileEvent",
            "ImRegistryEvent"
          ],
          "message": "`ActorUsername` is not a column of ImWebSession; it does exist on ImProcessCreate, ImAuthentication, ImFileEvent, ImRegistryEvent."
        },
        {
          "kind": "wrong_table",
          "field": "ActingProcessCommandLine",
          "tables_in_scope": [
            "ImWebSession"
          ],
          "exists_on": [
            "ImProcessCreate",
            "ImFileEvent",
            "ImRegistryEvent"
          ],
          "message": "`ActingProcessCommandLine` is not a column of ImWebSession; it does exist on ImProcessCreate, ImFileEvent, ImRegistryEvent."
        }
      ]
    }
  ]
}