# intel/.llm_actor_uc_cache/f2/f2261b530a93be699554a31ae09db5db084defca.json
{
  "name": "TeamTNT",
  "ucs": [
    {
      "title": "[LLM] TeamTNT cloud-credential harvest from compromised Linux/container hosts (IMDS + ~/.aws + kubelet creds)",
      "description": "Hunts TeamTNT's signature post-compromise behavior: a process running inside a container or on a Linux host enumerating AWS Instance Metadata Service (169.254.169.254), reading ~/.aws/credentials, ~/.docker/config.json, or /var/lib/kubelet tokens. TeamTNT (a.k.a. the same cluster Cado tracks as the 'Chimaera' / 'Black-T' / 'Hildegard' campaigns) uniquely chains cloud cred theft with XMRig rather than data theft, which separates it from generic IMDS abuse by red teams.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "actions",
      "conf": "high",
      "techs": [
        "T1552.001",
        "T1552.005",
        "T1552.004",
        "T1613",
        "T1611"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.user) as user from datamodel=Endpoint.Processes where Processes.os=Linux AND (Processes.process=\"*169.254.169.254*\" OR Processes.process=\"*/.aws/credentials*\" OR Processes.process=\"*/.aws/config*\" OR Processes.process=\"*/var/lib/kubelet*\" OR Processes.process=\"*/.docker/config.json*\" OR Processes.process=\"*AWS_ACCESS_KEY*\" OR Processes.process=\"*latest/meta-data/iam/security-credentials*\" OR Processes.process=\"*GOOGLE_APPLICATION_CREDENTIALS*\") by host Processes.process_name Processes.parent_process_name | `drop_dm_object_name(Processes)` | where parent_process_name IN (\"sh\",\"bash\",\"dash\",\"zsh\",\"ash\",\"docker\",\"dockerd\",\"containerd-shim\",\"containerd-shim-runc-v2\",\"runc\",\"containerd\",\"kubelet\",\"crio\",\"python\",\"python3\",\"perl\") | join type=outer host [| tstats summariesonly=true values(All_Traffic.dest_ip) as dest_ips count as net_count from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest_ip=\"169.254.169.254\" by All_Traffic.src as host] | eval risk=case(net_count>0 AND match(cmdline,\"\\.aws/credentials\"),\"high\", net_count>0,\"medium\", 1=1,\"low\") | table firstTime lastTime host user parent process_name cmdline dest_ips risk",
      "kql": "let credPaths = dynamic([\"/.aws/credentials\",\"/.aws/config\",\"/.docker/config.json\",\"/var/lib/kubelet\",\"/root/.kube/config\",\"/.gcloud/credentials\",\"/.azure/accessTokens.json\",\"id_rsa\",\"id_ecdsa\"]);\nlet credHunters = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where DeviceInfo_OS == \"Linux\" or InitiatingProcessFileName in~ (\"docker\",\"dockerd\",\"containerd-shim\",\"runc\",\"kubelet\",\"crio\",\"bash\",\"sh\",\"dash\",\"zsh\")\n| where ProcessCommandLine has_any (credPaths)\n     or ProcessCommandLine has \"169.254.169.254\"\n     or ProcessCommandLine has \"latest/meta-data/iam/security-credentials\"\n     or ProcessCommandLine matches regex @\"AWS_(ACCESS_KEY|SECRET_ACCESS_KEY|SESSION_TOKEN)\"\n| extend imds = ProcessCommandLine has \"169.254.169.254\", credFile = ProcessCommandLine has_any (credPaths)\n| project Timestamp, DeviceName, AccountName, InitiatingProcessFileName, FileName, ProcessCommandLine, imds, credFile, ReportId;\nlet imdsNet = DeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where RemoteIP == \"169.254.169.254\"\n| project NetTime=Timestamp, DeviceName, NetInitiatingProcess=InitiatingProcessFileName, NetCmd=InitiatingProcessCommandLine;\ncredHunters\n| join kind=leftouter (imdsNet) on DeviceName\n| where credFile or imds or isnotempty(NetTime)\n| extend Severity = case(credFile and (imds or isnotempty(NetTime)), \"High\", credFile or isnotempty(NetTime), \"Medium\", \"Low\")\n| order by Timestamp desc",
      "rationale": "TeamTNT is the only major actor whose entire economic model is built on harvesting cloud creds from container hosts to spread XMRig \u2014 Palo Alto Unit 42 (Black-T, Hildegard), Cado, Lacework, and Trend Micro reports all document the exact tool family (TNTbotinger, ezuri-packed shell scripts, Chimaera worm) calling IMDS and grepping ~/.aws/credentials. Joining IMDS network hits against credential-fil"
    },
    {
      "title": "[LLM] TeamTNT container escape \u2192 XMRig / Tsunami deployment chain (Docker/K8s API abuse + miner)",
      "description": "Hunts the canonical TeamTNT kill chain: an exposed Docker daemon or kubelet API spawns a privileged/host-mounted container, which curls a stage-2 shell script (often Ezuri-packed), disables defenses, and drops XMRig (or kdevtmpfsi/kinsing collisions, Tsunami IRC bot, Diamorphine rootkit) with persistence via systemd. This combination \u2014 container-runtime parent + miner binary + chattr/history-wipe + mining-pool egress \u2014 is a TeamTNT fingerprint.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1610",
        "T1611",
        "T1496.001",
        "T1543.002",
        "T1070.003",
        "T1222.002",
        "T1059.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.process_name) as proc values(Processes.parent_process_name) as parent from datamodel=Endpoint.Processes where Processes.os=Linux AND Processes.parent_process_name IN (\"dockerd\",\"docker\",\"containerd\",\"containerd-shim\",\"containerd-shim-runc-v2\",\"runc\",\"crio\",\"kubelet\") AND (Processes.process_name IN (\"xmrig\",\"kdevtmpfsi\",\"kinsing\",\"tsh\",\"masscan\",\"pnscan\",\"zgrab\",\"tmate\",\"weave\") OR Processes.process=\"*--donate-level*\" OR Processes.process=\"*stratum+tcp*\" OR Processes.process=\"*moneroocean*\" OR Processes.process=\"*supportxmr*\" OR Processes.process=\"*c3pool*\" OR Processes.process=\"*nanopool*\" OR Processes.process=\"*chattr -i*\" OR Processes.process=\"*history -c*\" OR Processes.process=\"*/dev/null 2>&1*xmr*\" OR Processes.process=\"*ezuri*\" OR Processes.process=\"*Diamorphine*\" OR Processes.process=\"*tnt*.sh*\" OR Processes.process=\"*kube.lateral.sh*\") by host Processes.user Processes.process_name Processes.parent_process_name | `drop_dm_object_name(Processes)` | append [| tstats summariesonly=true count from datamodel=Network_Traffic.All_Traffic where All_Traffic.dest_port IN (3333,4444,5555,7777,8888,9999,14444,14433,45700) by All_Traffic.src All_Traffic.dest All_Traffic.dest_port | `drop_dm_object_name(All_Traffic)` | rename src as host] | stats values(*) as * count by host | where count>=2",
      "kql": "let runtimeParents = dynamic([\"dockerd\",\"docker\",\"containerd\",\"containerd-shim\",\"containerd-shim-runc-v2\",\"runc\",\"crio\",\"kubelet\"]);\nlet minerNames = dynamic([\"xmrig\",\"kdevtmpfsi\",\"kinsing\",\"tsh\",\"masscan\",\"pnscan\",\"zgrab\",\"tmate\"]);\nlet poolStrings = dynamic([\"stratum+tcp\",\"--donate-level\",\"supportxmr\",\"moneroocean\",\"c3pool\",\"nanopool\",\"minexmr\",\"xmrpool\"]);\nlet defenseEvasion = dynamic([\"chattr -i\",\"chattr +i\",\"history -c\",\"unset HISTFILE\",\"ld.so.preload\",\"Diamorphine\",\"ezuri\"]);\nlet stageOne = DeviceProcessEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (runtimeParents) or InitiatingProcessParentFileName in~ (runtimeParents)\n| where FileName in~ (minerNames)\n   or ProcessCommandLine has_any (poolStrings)\n   or ProcessCommandLine has_any (defenseEvasion)\n   or ProcessCommandLine matches regex @\"(?i)(curl|wget)\\s+[^|]*\\|\\s*(sh|bash)\"\n| extend signal = case(\n    FileName in~ (minerNames), \"miner_binary\",\n    ProcessCommandLine has_any (poolStrings), \"mining_pool_arg\",\n    ProcessCommandLine has_any (defenseEvasion), \"defense_evasion\",\n    \"curl_pipe_sh\")\n| project Timestamp, DeviceName, AccountName, InitiatingProcessFileName, InitiatingProcessParentFileName, FileName, ProcessCommandLine, signal, ReportId;\nlet poolEgress = DeviceNetworkEvents\n| where Timestamp > ago(7d)\n| where RemotePort in (3333,4444,5555,7777,8888,9999,14444,14433,45700)\n   or RemoteUrl has_any (\"supportxmr\",\"moneroocean\",\"minexmr\",\"c3pool\",\"nanopool\",\"xmrpool\")\n| project NetTime=Timestamp, DeviceName, RemoteIP, RemotePort, RemoteUrl, InitiatingProcessFileName;\nstageOne\n| join kind=leftouter (poolEgress) on DeviceName\n| summarize signals=make_set(signal), pools=make_set(RemotePort), urls=make_set(RemoteUrl), cmds=make_set(ProcessCommandLine, 5), firstSeen=min(Timestamp), lastSeen=max(Timestamp) by DeviceName, InitiatingProcessFileName\n| where array_length(signals) >= 2 or array_length(pools) > 0\n| extend Severity = iff(array_length(pools)>0 and array_length(signals)>=2, \"Critical\", \"High\")",
      "rationale": "TeamTNT campaigns documented by Trend Micro (Chimaera), Lacework (TNTbotinger), AT&T Alien Labs (Chimaera), and Cado (the original 2020 worm) all share the same fingerprint: a process whose parent is a container runtime (not init/systemd) launching XMRig/Tsunami plus chattr/history wipes plus stratum egress. Requiring two of {miner binary, pool string, defense-evasion command, curl|sh} from a cont"
    }
  ]
}