# intel/.llm_actor_uc_cache/b6/b6aa26bd7d8f6754f249ab0ed7f4a0868f1bba71.json
{
  "name": "APT34",
  "ucs": [
    {
      "title": "[LLM] APT34 (OilRig) DNS-tunnel C2 \u2014 Helminth / DNSpionage / BondUpdater pattern",
      "description": "Hunts the high-entropy, long-subdomain DNS query bursts that APT34's signature implants (Helminth, ISMAgent, ALMA Communicator, DNSpionage, BondUpdater, Karkoff) use as their primary C2 channel. Every major OilRig campaign since 2016 has carried a DNS-tunnelling component, often as the *only* egress from segmented victim networks, which makes this one of the highest-fidelity behavioural pivots for this specific actor.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "c2",
      "conf": "high",
      "techs": [
        "T1071.004",
        "T1041",
        "T1090",
        "T1008"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count from datamodel=Network_Resolution where Network_Resolution.record_type IN (\"A\",\"AAAA\",\"TXT\",\"CNAME\") by _time span=10m Network_Resolution.src Network_Resolution.query Network_Resolution.record_type\n| rename Network_Resolution.* AS *\n| eval qlen=len(query), labels=split(query,\".\"), label_count=mvcount(labels), first_label=mvindex(labels,0), first_label_len=len(first_label)\n| where label_count >= 3 AND first_label_len >= 30\n| eval parent_domain=mvindex(labels,label_count-2).\".\".mvindex(labels,label_count-1)\n| stats dc(query) AS unique_subdomains, avg(first_label_len) AS avg_first_label_len, max(first_label_len) AS max_first_label_len, values(record_type) AS rtypes, values(query) AS sample_queries, count AS query_count by src parent_domain\n| where unique_subdomains >= 25 AND avg_first_label_len >= 30\n| sort - unique_subdomains\n| `security_content_ctime(_time)`",
      "kql": "// APT34 / OilRig \u2014 DNS tunnelling C2 (Helminth / DNSpionage / BondUpdater)\nlet _browsers = dynamic([\"chrome.exe\",\"msedge.exe\",\"firefox.exe\",\"iexplore.exe\",\"brave.exe\",\"opera.exe\",\"safari.exe\"]);\nlet _benign_resolvers = dynamic([\"svchost.exe\",\"searchindexer.exe\",\"backgroundtaskhost.exe\",\"smartscreen.exe\"]);\nDeviceEvents\n| where Timestamp > ago(7d)\n| where ActionType == \"DnsQueryResponse\"\n| extend QueryName = tolower(tostring(parse_json(AdditionalFields).QueryName))\n| where isnotempty(QueryName)\n| extend Labels = split(QueryName, \".\")\n| extend LabelCount = array_length(Labels)\n| where LabelCount >= 3                                  // skip TLD-only / 2-label\n| extend FirstLabel = tostring(Labels[0])\n| extend FirstLabelLen = strlen(FirstLabel)\n| where FirstLabelLen >= 30                              // OilRig encoded subdomains run 30-63 chars\n| extend ParentDomain = strcat(tostring(Labels[LabelCount-2]), \".\", tostring(Labels[LabelCount-1]))\n| where InitiatingProcessFileName !in~ (_browsers)       // browsers do legit long DNS\n| where InitiatingProcessFileName !in~ (_benign_resolvers)\n| where ParentDomain !endswith \"in-addr.arpa\"           // PTR lookups\n| where ParentDomain !has_any (\"akamai\",\"cloudfront\",\"azureedge\",\"trafficmanager\",\"office.com\",\"microsoft.com\",\"office365.com\",\"windows.com\",\"google.com\",\"googleapis.com\",\"gstatic.com\",\"apple.com\",\"icloud.com\")\n| summarize UniqueSubdomains = dcount(QueryName),\n            AvgFirstLabelLen = avg(FirstLabelLen),\n            MaxFirstLabelLen = max(FirstLabelLen),\n            QueryCount = count(),\n            Processes = make_set(InitiatingProcessFileName, 8),\n            ProcessCmds = make_set(InitiatingProcessCommandLine, 8),\n            SampleQueries = make_set(QueryName, 8)\n            by DeviceName, DeviceId, ParentDomain, bin(Timestamp, 1h)\n| where UniqueSubdomains >= 20 and AvgFirstLabelLen >= 30\n| order by UniqueSubdomains desc",
      "sentinel_kql": "// APT34 / OilRig \u2014 DNS tunnelling C2 via Windows DNS server logs\nlet _benign_parents = dynamic([\"akamai\",\"cloudfront\",\"azureedge\",\"trafficmanager\",\"office.com\",\"microsoft.com\",\"office365.com\",\"windows.com\",\"google.com\",\"googleapis.com\",\"gstatic.com\",\"apple.com\",\"icloud.com\"]);\nDnsEvents\n| where TimeGenerated > ago(7d)\n| where QueryType in (\"A\",\"AAAA\",\"TXT\",\"CNAME\")\n| extend QueryName = tolower(Name)\n| where isnotempty(QueryName)\n| extend Labels = split(QueryName, \".\")\n| extend LabelCount = array_length(Labels)\n| where LabelCount >= 3\n| extend FirstLabel = tostring(Labels[0])\n| extend FirstLabelLen = strlen(FirstLabel)\n| where FirstLabelLen >= 30\n| extend ParentDomain = strcat(tostring(Labels[LabelCount-2]), \".\", tostring(Labels[LabelCount-1]))\n| where ParentDomain !endswith \"in-addr.arpa\"\n| where ParentDomain !has_any (_benign_parents)\n| summarize UniqueSubdomains = dcount(QueryName),\n            AvgFirstLabelLen = avg(FirstLabelLen),\n            MaxFirstLabelLen = max(FirstLabelLen),\n            QueryCount = count(),\n            SampleQueries = make_set(QueryName, 8)\n            by ClientIP, Computer, ParentDomain, bin(TimeGenerated, 1h)\n| where UniqueSubdomains >= 20 and AvgFirstLabelLen >= 30\n| order by UniqueSubdomains desc",
      "sigma_yaml": "",
      "rationale": "DNS-tunnel C2 is APT34's *signature* egress channel \u2014 Palo Alto Unit 42, Cisco Talos and FireEye/Mandiant have all documented Helminth, DNSpionage, ALMA Communicator, BondUpdater and Karkoff using base32/hex-encoded subdomains 30-63 chars long against attacker-registered second-level domains. The combined predicate (long first label + non-browser/non-system parent + bursty unique-subdomain volume "
    },
    {
      "title": "[LLM] APT34 (OilRig) Office-macro \u2192 PowerShell/VBS stager \u2192 schtasks persistence chain (BondUpdater / Helminth initial access)",
      "description": "Hunts APT34's documented spearphishing chain where a weaponised .doc/.xls macro spawns cscript/wscript or powershell, which then registers a scheduled task for backdoor persistence within minutes. This exact 3-step sequence (Office \u2192 script-host/PowerShell \u2192 schtasks) is the hallmark Helminth, BondUpdater and QUADAGENT delivery pattern reported by FireEye, Cisco Talos, ClearSky and Microsoft (Hazel Sandstorm).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1566.001",
        "T1059.001",
        "T1059.005",
        "T1059.007",
        "T1053.005",
        "T1204.002"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count min(_time) AS stager_time values(Processes.process) AS stager_cmd values(Processes.parent_process_name) AS office_parent from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"eqnedt32.exe\") Processes.process_name IN (\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"cscript.exe\",\"wscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\") by Processes.dest Processes.user Processes.process_name\n| `drop_dm_object_name(Processes)`\n| join type=inner dest user [\n    | tstats `security_content_summariesonly` min(_time) AS schtask_time values(Processes.process) AS schtask_cmd from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" (Processes.process=\"*create*\" OR Processes.process=\"*CREATE*\" OR Processes.process=\"*/sc *\") by Processes.dest Processes.user\n    | `drop_dm_object_name(Processes)`\n  ]\n| eval delta_sec = abs(schtask_time - stager_time)\n| where delta_sec <= 600\n| table dest user office_parent process_name stager_cmd schtask_cmd delta_sec\n| `security_content_ctime(stager_time)`",
      "kql": "// APT34 / OilRig \u2014 Office macro stager + scheduled-task persistence inside 10 minutes\nlet LookbackDays = 7d;\nlet ChainWindow = 10m;\nlet OfficeParents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"eqnedt32.exe\"]);\nlet Stagers = dynamic([\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"cscript.exe\",\"wscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\"]);\nlet MacroSpawn =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where InitiatingProcessFileName in~ (OfficeParents)\n    | where FileName in~ (Stagers)\n    | where AccountName !endswith \"$\"\n    | project StagerTime = Timestamp, DeviceId, DeviceName, AccountName,\n              OfficeParent = InitiatingProcessFileName,\n              OfficeCmd = InitiatingProcessCommandLine,\n              StagerImage = FileName,\n              StagerCmd = ProcessCommandLine,\n              StagerSHA256 = SHA256;\nDeviceProcessEvents\n| where Timestamp > ago(LookbackDays)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has_any (\"/create\",\"/CREATE\",\"-create\",\"/sc \",\"/SC \")\n| where AccountName !endswith \"$\"\n| project SchTime = Timestamp, DeviceId, AccountName, SchCmd = ProcessCommandLine,\n          SchParent = InitiatingProcessFileName, SchParentCmd = InitiatingProcessCommandLine\n| join kind=inner MacroSpawn on DeviceId, AccountName\n| where SchTime between (StagerTime .. StagerTime + ChainWindow)\n| extend ChainDelaySec = datetime_diff('second', SchTime, StagerTime)\n| project StagerTime, SchTime, ChainDelaySec,\n          DeviceName, AccountName,\n          OfficeParent, OfficeCmd,\n          StagerImage, StagerCmd, StagerSHA256,\n          SchParent, SchParentCmd, SchCmd\n| order by StagerTime desc",
      "sentinel_kql": "// APT34 / OilRig \u2014 Office macro stager + schtasks persistence (Sentinel ASIM, falls back to SecurityEvent 4688)\nlet LookbackDays = 7d;\nlet ChainWindow = 10m;\nlet OfficeParents = dynamic([\"winword.exe\",\"excel.exe\",\"powerpnt.exe\",\"outlook.exe\",\"eqnedt32.exe\"]);\nlet Stagers = dynamic([\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\",\"cscript.exe\",\"wscript.exe\",\"mshta.exe\",\"rundll32.exe\",\"regsvr32.exe\"]);\nlet MacroSpawn =\n    ImProcessCreate\n    | where TimeGenerated > ago(LookbackDays)\n    | where ParentProcessName has_any (OfficeParents)\n    | where TargetProcessName has_any (Stagers)\n    | where ActorUsername !endswith \"$\"\n    | extend StagerImageName = tolower(tostring(split(TargetProcessName, \"\\\\\")[-1]))\n    | extend OfficeParentName = tolower(tostring(split(ParentProcessName, \"\\\\\")[-1]))\n    | project StagerTime = TimeGenerated, DvcHostname, ActorUsername,\n              OfficeParentName, ParentProcessCommandLine,\n              StagerImageName, TargetProcessCommandLine;\nImProcessCreate\n| where TimeGenerated > ago(LookbackDays)\n| where TargetProcessName endswith \"schtasks.exe\"\n| where TargetProcessCommandLine has_any (\"/create\",\"/CREATE\",\"-create\",\"/sc \",\"/SC \")\n| where ActorUsername !endswith \"$\"\n| project SchTime = TimeGenerated, DvcHostname, ActorUsername,\n          SchCmd = TargetProcessCommandLine,\n          SchParent = ParentProcessName\n| join kind=inner MacroSpawn on DvcHostname, ActorUsername\n| where SchTime between (StagerTime .. StagerTime + ChainWindow)\n| extend ChainDelaySec = datetime_diff('second', SchTime, StagerTime)\n| project StagerTime, SchTime, ChainDelaySec,\n          DvcHostname, ActorUsername,\n          OfficeParentName, ParentProcessCommandLine,\n          StagerImageName, TargetProcessCommandLine,\n          SchParent, SchCmd\n| order by StagerTime desc",
      "sigma_yaml": "title: APT34 OilRig - Office Application Spawning Script Host or PowerShell Stager\nid: 5b3b8e30-2f2c-4c4a-8b35-1e7a6f3c4d12\nstatus: experimental\ndescription: Detects Microsoft Office (Word/Excel/PowerPoint/Outlook) or Equation Editor spawning a scripting / LOLBin stager - the first hop of APT34 Helminth/BondUpdater/QUADAGENT delivery.\nreferences:\n  - https://attack.mitre.org/groups/G0049/\n  - https://unit42.paloaltonetworks.com/dns-tunneling-in-the-wild-overview-of-oilrigs-dns-tunneling/\nauthor: SOC Detection Engineering\ndate: 2026/05/06\ntags:\n  - attack.initial_access\n  - attack.execution\n  - attack.t1566.001\n  - attack.t1059.001\n  - attack.t1059.005\n  - attack.g0049\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_parent:\n    ParentImage|endswith:\n      - '\\winword.exe'\n      - '\\excel.exe'\n      - '\\powerpnt.exe'\n      - '\\outlook.exe'\n      - '\\eqnedt32.exe'\n  selection_child:\n    Image|endswith:\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\cmd.exe'\n      - '\\cscript.exe'\n      - '\\wscript.exe'\n      - '\\mshta.exe'\n      - '\\rundll32.exe'\n      - '\\regsvr32.exe'\n  filter_machine_acct:\n    User|endswith: '$'\n  condition: selection_parent and selection_child and not filter_machine_acct\nfalsepositives:\n  - Legitimate enterprise Office add-ins that shell out to script hosts\n  - Internal macro-driven automation (scope to user populations that should never run macros)\nlevel: high",
      "rationale": "Public technical write-ups on QUADAGENT (FireEye), BondUpdater (Palo Alto Unit 42), Helminth (Cisco Talos) and OopsIE (ClearSky) all describe the *same* three-step delivery \u2014 phished Office doc \u2192 cscript/wscript or powershell stager \u2192 schtasks-registered persistence task. Joining the macro-spawn event to a same-host, same-user schtasks /create within a 10-minute window encodes that exact chain and"
    }
  ]
}