# intel/.llm_actor_uc_cache/5d/5dd2354c538365526b217b55bcc8f3d4894a1d2a.json
{
  "name": "Lazarus Group",
  "ucs": [
    {
      "title": "[LLM] Lazarus / Diamond Sleet DLL search-order hijack via trusted Microsoft-signed binary in user-writable path",
      "description": "Hunts Lazarus Group's signature 'sideload from a writable folder' tradecraft \u2014 dropping a legitimate Microsoft / OEM-signed host binary (wsmprovhost.exe Comebacker loader, credwiz.exe DUSER sideload, iviewers.exe Manuscrypt loader, colorcpl.exe SmoothShifter, GoogleUpdate.exe / goopdate.dll Diamond Sleet pair, vmtoolsd.exe, snmptrap.exe, GUP.exe) into %APPDATA%, %TEMP%, %PUBLIC% or %PROGRAMDATA% alongside a malicious DLL. The Microsoft-signed loader gets hash-trust and AppLocker bypass; the DLL holds the real payload. This pair is the recurring backbone of Operation Dream Job, AppleJeus, the J",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.001",
        "T1574.002",
        "T1036.005",
        "T1027.013"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.process_hash) as host_binary_hash from datamodel=Endpoint.Processes where Processes.process_name IN (\"wsmprovhost.exe\",\"credwiz.exe\",\"iviewers.exe\",\"colorcpl.exe\",\"DTAUDIT.exe\",\"vmtoolsd.exe\",\"GoogleUpdate.exe\",\"GUP.exe\",\"snmptrap.exe\",\"RTKauUService64.exe\") AND (Processes.process_path=\"*\\\\Users\\\\Public\\\\*\" OR Processes.process_path=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Processes.process_path=\"*\\\\AppData\\\\Roaming\\\\*\" OR Processes.process_path=\"*\\\\ProgramData\\\\*\" OR Processes.process_path=\"*\\\\Windows\\\\Temp\\\\*\" OR Processes.process_path=\"*\\\\PerfLogs\\\\*\" OR Processes.process_path=\"*\\\\Users\\\\*\\\\Downloads\\\\*\") AND Processes.parent_process_name!=\"msiexec.exe\" AND Processes.parent_process_name!=\"setup.exe\" by Processes.dest Processes.user Processes.process_name Processes.process_path | `drop_dm_object_name(Processes)` | rex field=process_path \"(?<host_dir>.*)\\\\\\\\[^\\\\\\\\]+$\" | join type=inner dest host_dir [ | tstats summariesonly=true values(Filesystem.file_hash) as dll_hash from datamodel=Endpoint.Filesystem where Filesystem.file_name=\"*.dll\" by Filesystem.dest Filesystem.file_path Filesystem.file_name | `drop_dm_object_name(Filesystem)` | rex field=file_path \"(?<host_dir>.*)\\\\\\\\[^\\\\\\\\]+$\" ] | where host_binary_hash!=dll_hash | convert ctime(firstTime) ctime(lastTime) | table firstTime lastTime dest user process_name process_path file_name file_path host_binary_hash dll_hash parent cmdline",
      "kql": "let LazarusSideloadHosts = dynamic([\"wsmprovhost.exe\",\"credwiz.exe\",\"iviewers.exe\",\"colorcpl.exe\",\"DTAUDIT.exe\",\"vmtoolsd.exe\",\"GoogleUpdate.exe\",\"GUP.exe\",\"snmptrap.exe\",\"RTKauUService64.exe\"]);\nlet WritableStaging = dynamic([@\"\\Users\\Public\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\AppData\\Roaming\\\", @\"\\ProgramData\\\", @\"\\Windows\\Temp\\\", @\"\\PerfLogs\\\", @\"\\Downloads\\\"]);\nDeviceImageLoadEvents\n| where Timestamp > ago(7d)\n| where InitiatingProcessFileName in~ (LazarusSideloadHosts)\n| where InitiatingProcessFolderPath has_any (WritableStaging)\n// adjacent DLL \u2014 same folder as the trusted host binary\n| where FolderPath =~ InitiatingProcessFolderPath\n| where FileName endswith \".dll\"\n| where SHA256 != InitiatingProcessSHA256\n// strip well-known side-by-side runtime libs that are commonly bundled with installers\n| where FileName !in~ (\"dbghelp.dll\",\"msvcr100.dll\",\"msvcp100.dll\",\"vcruntime140.dll\",\"vcruntime140_1.dll\",\"msvcp140.dll\",\"ucrtbase.dll\")\n| join kind=leftouter (\n    DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | project DeviceId, ProcessId = InitiatingProcessId, ProcessCmd = InitiatingProcessCommandLine, ParentBinary = InitiatingProcessParentFileName\n  ) on $left.InitiatingProcessId == $right.ProcessId, DeviceId\n| project Timestamp, DeviceName, InitiatingProcessAccountName,\n          HostBinary = InitiatingProcessFileName,\n          HostPath = InitiatingProcessFolderPath,\n          HostSHA256 = InitiatingProcessSHA256,\n          SideloadedDll = FileName,\n          SideloadedDllSHA256 = SHA256,\n          ParentBinary, ProcessCmd\n| order by Timestamp desc",
      "sentinel_kql": "let LazarusSideloadHosts = dynamic([\"wsmprovhost.exe\",\"credwiz.exe\",\"iviewers.exe\",\"colorcpl.exe\",\"DTAUDIT.exe\",\"vmtoolsd.exe\",\"GoogleUpdate.exe\",\"GUP.exe\",\"snmptrap.exe\",\"RTKauUService64.exe\"]);\nlet WritableStaging = dynamic([@\"\\Users\\Public\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\AppData\\Roaming\\\", @\"\\ProgramData\\\", @\"\\Windows\\Temp\\\", @\"\\PerfLogs\\\", @\"\\Downloads\\\"]);\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| extend BinaryName = tostring(split(NewProcessName, \"\\\\\")[-1])\n| where tolower(BinaryName) in (LazarusSideloadHosts)\n| where NewProcessName has_any (WritableStaging)\n| where ParentProcessName !endswith \"\\\\msiexec.exe\"\n   and ParentProcessName !endswith \"\\\\setup.exe\"\n   and ParentProcessName !endswith \"\\\\TrustedInstaller.exe\"\n| project TimeGenerated, Computer, Account, BinaryName, NewProcessName, CommandLine, ParentProcessName, TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Lazarus DLL Search-Order Hijack via Trusted Microsoft-Signed Binary in User-Writable Path\nid: 7ad3d2e0-6a9f-4f1b-9e1d-laz001sideload\nstatus: experimental\ndescription: Detects DPRK Lazarus Group / Diamond Sleet DLL sideloading by spotting known-abused signed host binaries (wsmprovhost.exe, credwiz.exe, iviewers.exe, colorcpl.exe, GoogleUpdate.exe, vmtoolsd.exe, GUP.exe, snmptrap.exe) executing from user-writable directories instead of their canonical install paths. Used in Operation Dream Job, AppleJeus, the JumpCloud intrusion, and the TeamCity CVE-2023-42793 chain.\nreferences:\n  - https://attack.mitre.org/groups/G0032/\n  - https://www.microsoft.com/en-us/security/blog/2023/10/18/multiple-north-korean-threat-actors-exploiting-the-teamcity-cve-2023-42793-vulnerability/\n  - https://www.mandiant.com/resources/blog/lazarus-attacks-jumpcloud\nauthor: Detection Engineering\ndate: 2026/05/10\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_image:\n    Image|endswith:\n      - '\\wsmprovhost.exe'\n      - '\\credwiz.exe'\n      - '\\iviewers.exe'\n      - '\\colorcpl.exe'\n      - '\\DTAUDIT.exe'\n      - '\\vmtoolsd.exe'\n      - '\\GoogleUpdate.exe'\n      - '\\GUP.exe'\n      - '\\snmptrap.exe'\n      - '\\RTKauUService64.exe'\n  selection_path:\n    Image|contains:\n      - '\\Users\\Public\\'\n      - '\\AppData\\Local\\Temp\\'\n      - '\\AppData\\Roaming\\'\n      - '\\ProgramData\\'\n      - '\\Windows\\Temp\\'\n      - '\\Downloads\\'\n      - '\\PerfLogs\\'\n  filter_installer:\n    ParentImage|endswith:\n      - '\\msiexec.exe'\n      - '\\setup.exe'\n      - '\\TrustedInstaller.exe'\n  condition: selection_image and selection_path and not filter_installer\nfalsepositives:\n  - Legitimate developer staging or installer self-extraction prior to copy-to-Program-Files\n  - Portable software intentionally run from user folders (rare for these specific binaries)\nlevel: high\ntags:\n  - attack.defense_evasion\n  - attack.t1574.001\n  - attack.t1574.002\n  - attack.t1036.005\n  - attack.g0032",
      "rationale": "Lazarus / Diamond Sleet has used this exact 'trusted-loader-in-writable-path + adjacent malicious DLL' construct continuously since 2020 \u2014 Comebacker (wsmprovhost), DUSER (credwiz), Manuscrypt (iviewers), SmoothShifter (colorcpl), and the 2023 TeamCity / JumpCloud intrusions. The detection deliberately keys on the unusual launch location of legitimate signed binaries \u2014 not on the DLL itself, since",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] Lazarus FudModule BYOVD \u2014 vulnerable signed kernel driver service-install for admin-to-kernel EDR kill",
      "description": "Hunts Lazarus's FudModule rootkit installation pattern: registering a Windows service whose ImagePath points at a known-vulnerable signed kernel driver (Dell dbutil_2_3.sys, ENE ene.sys, MSI Afterburner RTCore64.sys, OpenLibSys WinRing0x64.sys, Intel iqvw64e.sys, Gigabyte gdrv.sys, ASRock AsrDrv101.sys, AppID downgrade) \u2014 used by Lazarus / Diamond Sleet across 2022-2024 to defeat Defender / Kaspersky / AhnLab kernel callbacks via the FudModule rootkit (CVE-2024-21338 AppLocker driver IOCTL exploitation, plus the older BYOVD chain).",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1543.003",
        "T1068",
        "T1014",
        "T1562.001"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Registry.process_name) as installer values(Registry.user) as user values(Registry.registry_value_data) as image_path from datamodel=Endpoint.Registry where Registry.registry_path=\"*\\\\System\\\\CurrentControlSet\\\\Services\\\\*\\\\ImagePath\" AND (Registry.registry_value_data=\"*dbutil_2_3.sys*\" OR Registry.registry_value_data=\"*DBUtil_2_3.sys*\" OR Registry.registry_value_data=\"*ene.sys*\" OR Registry.registry_value_data=\"*RTCore64.sys*\" OR Registry.registry_value_data=\"*WinRing0x64.sys*\" OR Registry.registry_value_data=\"*iqvw64e.sys*\" OR Registry.registry_value_data=\"*gdrv.sys*\" OR Registry.registry_value_data=\"*MsIo64.sys*\" OR Registry.registry_value_data=\"*AsrDrv101.sys*\" OR Registry.registry_value_data=\"*procexp152.sys*\" OR Registry.registry_value_data=\"*viragt64.sys*\" OR Registry.registry_value_data=\"*amsdk.sys*\") by Registry.dest Registry.registry_path | `drop_dm_object_name(Registry)` | convert ctime(firstTime) ctime(lastTime) | table firstTime lastTime dest user installer registry_path image_path",
      "kql": "let FudModuleBYOVD = dynamic([\"dbutil_2_3.sys\",\"ene.sys\",\"rtcore64.sys\",\"winring0x64.sys\",\"iqvw64e.sys\",\"gdrv.sys\",\"msio64.sys\",\"asrdrv101.sys\",\"procexp152.sys\",\"viragt64.sys\",\"amsdk.sys\",\"appid.sys\"]);\nlet RegistryInstall =\n    DeviceRegistryEvents\n    | where Timestamp > ago(30d)\n    | where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n    | where RegistryKey startswith @\"HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Services\\\"\n    | where RegistryValueName =~ \"ImagePath\"\n    | extend DriverFile = tolower(tostring(split(RegistryValueData, \"\\\\\")[-1]))\n    | where DriverFile in (FudModuleBYOVD)\n    | where InitiatingProcessFileName !in~ (\"TrustedInstaller.exe\",\"msiexec.exe\")\n    | project Timestamp, DeviceName, ActionType, RegistryKey, DriverImagePath = RegistryValueData,\n              DriverFile, Installer = InitiatingProcessFileName,\n              InstallerCmd = InitiatingProcessCommandLine,\n              InstallerUser = InitiatingProcessAccountName, InitiatingProcessSHA256;\nlet DriverDrop =\n    DeviceFileEvents\n    | where Timestamp > ago(30d)\n    | where ActionType == \"FileCreated\"\n    | extend DriverFile = tolower(FileName)\n    | where DriverFile in (FudModuleBYOVD)\n    | where FolderPath !startswith @\"C:\\Windows\\System32\\drivers\\\"\n         or InitiatingProcessFolderPath has_any (@\"\\Users\\\", @\"\\ProgramData\\\", @\"\\AppData\\\", @\"\\Public\\\", @\"\\Temp\\\")\n    | project Timestamp, DeviceName, ActionType, FolderPath, FileName, DriverFile, SHA256,\n              Dropper = InitiatingProcessFileName, DropperCmd = InitiatingProcessCommandLine,\n              DropperUser = InitiatingProcessAccountName;\nunion RegistryInstall, DriverDrop\n| order by Timestamp desc",
      "sentinel_kql": "let FudModuleBYOVD = dynamic([\"dbutil_2_3.sys\",\"ene.sys\",\"rtcore64.sys\",\"winring0x64.sys\",\"iqvw64e.sys\",\"gdrv.sys\",\"msio64.sys\",\"asrdrv101.sys\",\"procexp152.sys\",\"viragt64.sys\",\"amsdk.sys\",\"appid.sys\"]);\n// Service install via Registry (System / Sysmon EID 13)\nlet RegInstall = WindowsEvent\n    | where TimeGenerated > ago(30d)\n    | where Provider == \"Microsoft-Windows-Sysmon\" and EventID == 13\n    | extend TargetObject = tostring(EventData.TargetObject),\n             Details = tostring(EventData.Details),\n             Image = tostring(EventData.Image)\n    | where TargetObject startswith @\"HKLM\\System\\CurrentControlSet\\Services\\\"\n         and TargetObject endswith @\"\\ImagePath\"\n    | extend DriverFile = tolower(tostring(split(Details, \"\\\\\")[-1]))\n    | where DriverFile in (FudModuleBYOVD)\n    | project TimeGenerated, Computer, Source = \"Sysmon Registry\", DriverFile, DriverImagePath = Details, Installer = Image, RegistryKey = TargetObject;\n// Backup via classic Windows service-install event (System log EID 7045)\nlet ServiceInstall = Event\n    | where TimeGenerated > ago(30d)\n    | where EventLog == \"System\" and EventID == 7045\n    | extend ImagePath = tostring(extract(@\"<Data Name=['\\\"]ImagePath['\\\"]>([^<]+)\", 1, EventData)),\n             ServiceName = tostring(extract(@\"<Data Name=['\\\"]ServiceName['\\\"]>([^<]+)\", 1, EventData))\n    | extend DriverFile = tolower(tostring(split(ImagePath, \"\\\\\")[-1]))\n    | where DriverFile in (FudModuleBYOVD)\n    | project TimeGenerated, Computer, Source = \"System 7045\", DriverFile, DriverImagePath = ImagePath, Installer = \"\", RegistryKey = ServiceName;\nunion RegInstall, ServiceInstall\n| order by TimeGenerated desc",
      "sigma_yaml": "title: Lazarus FudModule BYOVD - Vulnerable Kernel Driver Service Install\nid: 8b14e7a2-2c3f-44a6-8b1c-laz002fudmod\nstatus: experimental\ndescription: Detects DPRK Lazarus Group / Diamond Sleet installing a Windows kernel-driver service whose ImagePath points at a known-vulnerable signed driver used by the FudModule rootkit (Dell dbutil_2_3.sys, ENE ene.sys, MSI RTCore64.sys, OpenLibSys WinRing0x64.sys, Intel iqvw64e.sys, Gigabyte gdrv.sys, ASRock AsrDrv101.sys, etc.) for admin-to-kernel privilege escalation and EDR kernel-callback removal.\nreferences:\n  - https://attack.mitre.org/groups/G0032/\n  - https://www.microsoft.com/en-us/security/blog/2024/02/28/lazarus-fudmodule-rootkit-cve-2024-21338/\n  - https://decoded.avast.io/janvojtesek/lazarus-and-the-fudmodule-rootkit-beyond-byovd-with-an-admin-to-kernel-zero-day/\n  - https://www.welivesecurity.com/2022/09/30/amazon-themed-campaigns-lazarus-netherlands-belgium/\nauthor: Detection Engineering\ndate: 2026/05/10\nlogsource:\n  category: registry_set\n  product: windows\ndetection:\n  selection_key:\n    TargetObject|startswith: 'HKLM\\System\\CurrentControlSet\\Services\\'\n    TargetObject|endswith: '\\ImagePath'\n  selection_driver:\n    Details|contains:\n      - 'dbutil_2_3.sys'\n      - 'ene.sys'\n      - 'RTCore64.sys'\n      - 'WinRing0x64.sys'\n      - 'iqvw64e.sys'\n      - 'gdrv.sys'\n      - 'MsIo64.sys'\n      - 'AsrDrv101.sys'\n      - 'procexp152.sys'\n      - 'viragt64.sys'\n      - 'amsdk.sys'\n  filter_legit_installer:\n    Image|endswith:\n      - '\\TrustedInstaller.exe'\n      - '\\msiexec.exe'\n  condition: selection_key and selection_driver and not filter_legit_installer\nfalsepositives:\n  - Hardware-vendor utility (Dell Command Update, MSI Afterburner, Gigabyte Control Center) install via vendor MSI \u2014 verify Image is a signed vendor installer\n  - Legitimate IT troubleshooting (Process Explorer driver) on engineer endpoints\nlevel: critical\ntags:\n  - attack.defense_evasion\n  - attack.privilege_escalation\n  - attack.t1543.003\n  - attack.t1068\n  - attack.t1014\n  - attack.t1562.001\n  - attack.g0032",
      "rationale": "Microsoft (Feb 2024) and Avast (Dec 2023) attribute the FudModule rootkit exclusively to Lazarus / Diamond Sleet, with this exact BYOVD chain documented across the 3CX intrusion (Dell DBUtil), the Aug 2023 aerospace campaigns (ene.sys), and the Feb 2024 CVE-2024-21338 admin-to-kernel exploitation (vulnerable AppID.sys downgrade). Lazarus is one of the few state actors whose tradecraft REQUIRES ker",
      "_sentinel_field_issues": [
        {
          "kind": "unknown_field",
          "field": "System",
          "tables_in_scope": [
            "Event",
            "WindowsEvent"
          ],
          "suggestion": null,
          "message": "`System` is not a canonical column of Event, WindowsEvent."
        },
        {
          "kind": "wrong_table",
          "field": "ServiceName",
          "tables_in_scope": [
            "Event",
            "WindowsEvent"
          ],
          "exists_on": [
            "SecurityEvent"
          ],
          "message": "`ServiceName` is not a column of Event, WindowsEvent; it does exist on SecurityEvent."
        }
      ],
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    }
  ]
}