# intel/.llm_actor_uc_cache/af/af15559187fa5b2a5528510a6df7cc33d5dee923.json
{
  "name": "APT10",
  "ucs": [
    {
      "title": "[LLM] APT10 / Stone Panda DLL side-loading via legitimate signed binary in non-standard path (PlugX / ANEL / UPPERCUT loader)",
      "description": "Hunts APT10's signature delivery technique: dropping a legitimately-signed Microsoft / Java / McAfee / Google binary into a writable directory (ProgramData, AppData, Public) alongside a maliciously-replaced sibling DLL (jli.dll, msvcr100.dll, libcurl.dll, oci.dll, etc.) so the binary side-loads attacker code into a trusted process. This staging pattern is documented across PlugX, ANEL, ChChes, RedLeaves and UPPERCUT campaigns attributed to MenuPass / Stone Panda / Bronze Riverside.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1574.002",
        "T1027",
        "T1036.005",
        "T1105"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count min(_time) as firstTime max(_time) as lastTime from datamodel=Endpoint.Processes where Processes.process_name IN (\"jjs.exe\",\"javaw.exe\",\"GoogleUpdate.exe\",\"mcoemcpy.exe\",\"mcods.exe\",\"DPMonitor.exe\",\"FlashHelp_Player.exe\",\"PotPlayerMini.exe\",\"vmnat.exe\",\"hkcmd.exe\") AND (Processes.process_path=\"*\\\\ProgramData\\\\*\" OR Processes.process_path=\"*\\\\AppData\\\\*\" OR Processes.process_path=\"*\\\\Users\\\\Public\\\\*\" OR Processes.process_path=\"*\\\\PerfLogs\\\\*\" OR Processes.process_path=\"*\\\\Windows\\\\Temp\\\\*\" OR Processes.process_path=\"*\\\\Intel\\\\*\") by host Processes.dest Processes.user Processes.process_name Processes.process_path Processes.process_hash Processes.parent_process_name | `drop_dm_object_name(Processes)` | rex field=process_path \"(?<loader_dir>.+)\\\\\\\\[^\\\\\\\\]+$\" | join type=inner dest loader_dir [ | tstats `security_content_summariesonly` count from datamodel=Endpoint.Filesystem where Filesystem.file_name IN (\"jli.dll\",\"msvcr100.dll\",\"mscoree.dll\",\"libcurl.dll\",\"oci.dll\",\"secur32.dll\",\"cryptbase.dll\",\"version.dll\",\"wtsapi32.dll\",\"netapi32.dll\",\"goopdate.dll\",\"FlashHelp.dll\",\"PotPlayer.dll\",\"vmtools.dll\") Filesystem.action=created by Filesystem.dest Filesystem.file_path Filesystem.file_name Filesystem.file_hash | `drop_dm_object_name(Filesystem)` | rex field=file_path \"(?<loader_dir>.+)\\\\\\\\[^\\\\\\\\]+$\" ] | table firstTime lastTime dest user process_name process_path process_hash parent_process_name file_name file_hash | `security_content_ctime(firstTime)` | `security_content_ctime(lastTime)`",
      "kql": "// APT10 / Stone Panda \u2014 legit binary in writable dir side-loading a malicious sibling DLL\nlet StagingDirs = dynamic([@\"\\ProgramData\\\", @\"\\AppData\\Local\\\", @\"\\AppData\\Roaming\\\", @\"\\Users\\Public\\\", @\"\\PerfLogs\\\", @\"\\Windows\\Temp\\\", @\"\\Intel\\\"]);\nlet SideLoaders = dynamic([\"jjs.exe\",\"javaw.exe\",\"java.exe\",\"GoogleUpdate.exe\",\"mcoemcpy.exe\",\"mcods.exe\",\"DPMonitor.exe\",\"FlashHelp_Player.exe\",\"PotPlayerMini.exe\",\"vmnat.exe\",\"hkcmd.exe\",\"avp.exe\",\"OInfoP11.exe\"]);\nlet APT10Dlls = dynamic([\"jli.dll\",\"msvcr100.dll\",\"mscoree.dll\",\"libcurl.dll\",\"oci.dll\",\"secur32.dll\",\"cryptbase.dll\",\"version.dll\",\"wtsapi32.dll\",\"netapi32.dll\",\"goopdate.dll\",\"FlashHelp.dll\",\"PotPlayer.dll\",\"vmtools.dll\"]);\nDeviceImageLoadEvents\n| where Timestamp > ago(30d)\n| where FileName in~ (APT10Dlls)\n| where FolderPath has_any (StagingDirs)                            // DLL dropped into writable staging dir\n| where InitiatingProcessFileName in~ (SideLoaders)                 // loader is a known side-load target\n| where InitiatingProcessFolderPath !startswith @\"C:\\Windows\\\"     // not the genuine system copy\n   and InitiatingProcessFolderPath !startswith @\"C:\\Program Files\"\n| where tolower(InitiatingProcessFolderPath) == tolower(FolderPath) // adjacent \u2014 classic side-load layout\n| where InitiatingProcessAccountName !endswith \"$\"\n| project Timestamp, DeviceName, InitiatingProcessAccountName,\n          Loader     = InitiatingProcessFileName,\n          LoaderPath = InitiatingProcessFolderPath,\n          LoaderHash = InitiatingProcessSHA256,\n          LoaderCmd  = InitiatingProcessCommandLine,\n          DLL        = FileName,\n          DllPath    = FolderPath,\n          DllHash    = SHA256\n| order by Timestamp desc",
      "sentinel_kql": "// APT10 / Stone Panda \u2014 legit signed binary in writable dir side-loading malicious sibling DLL (Sysmon EID 7 image-load)\nlet StagingDirs = dynamic([@\"\\ProgramData\\\", @\"\\AppData\\Local\\\", @\"\\AppData\\Roaming\\\", @\"\\Users\\Public\\\", @\"\\PerfLogs\\\", @\"\\Windows\\Temp\\\", @\"\\Intel\\\"]);\nlet SideLoaders = dynamic([\"jjs.exe\",\"javaw.exe\",\"java.exe\",\"GoogleUpdate.exe\",\"mcoemcpy.exe\",\"mcods.exe\",\"DPMonitor.exe\",\"FlashHelp_Player.exe\",\"PotPlayerMini.exe\",\"vmnat.exe\",\"hkcmd.exe\"]);\nlet APT10Dlls = dynamic([\"jli.dll\",\"msvcr100.dll\",\"mscoree.dll\",\"libcurl.dll\",\"oci.dll\",\"secur32.dll\",\"cryptbase.dll\",\"version.dll\",\"wtsapi32.dll\",\"netapi32.dll\",\"goopdate.dll\",\"FlashHelp.dll\",\"PotPlayer.dll\",\"vmtools.dll\"]);\nWindowsEvent\n| where TimeGenerated > ago(30d)\n| where Provider == \"Microsoft-Windows-Sysmon\" and EventID == 7\n| extend Image       = tostring(EventData.Image),\n         ImageLoaded = tostring(EventData.ImageLoaded),\n         User        = tostring(EventData.User),\n         Signed      = tostring(EventData.Signed),\n         Hashes      = tostring(EventData.Hashes)\n| extend LoaderName = tolower(extract(@\"\\\\([^\\\\]+)$\", 1, Image)),\n         DllName    = tolower(extract(@\"\\\\([^\\\\]+)$\", 1, ImageLoaded)),\n         LoaderDir  = tolower(extract(@\"^(.+)\\\\[^\\\\]+$\", 1, Image)),\n         DllDir     = tolower(extract(@\"^(.+)\\\\[^\\\\]+$\", 1, ImageLoaded))\n| where DllName in~ (APT10Dlls)\n| where LoaderName in~ (SideLoaders)\n| where ImageLoaded has_any (StagingDirs)\n| where not(Image startswith @\"C:\\Windows\\\") and not(Image startswith @\"C:\\Program Files\")\n| where LoaderDir == DllDir\n| where User !endswith \"$\"\n| project TimeGenerated, Computer, User, LoaderName, Image, DllName, ImageLoaded, Signed, Hashes\n| order by TimeGenerated desc",
      "sigma_yaml": "title: APT10 Stone Panda DLL Side-Loading via Legitimate Binary in Writable Path\nid: a1c47e5e-7f2c-4d8b-9c1a-1f9c3a23d101\nstatus: experimental\ndescription: Detects APT10 / MenuPass / Stone Panda DLL side-loading where a legitimately-signed binary (Java jjs/javaw, GoogleUpdate, McAfee mcoemcpy, FlashHelp_Player, PotPlayer) executes from a non-standard writable path (ProgramData, AppData, Public, PerfLogs, Temp) and loads a sibling DLL with a name historically abused by ANEL, UPPERCUT, ChChes, RedLeaves and PlugX (jli.dll, msvcr100.dll, libcurl.dll, oci.dll, etc.).\nreferences:\n  - https://attack.mitre.org/groups/G0045/\n  - https://www.pwc.co.uk/cyber-security/assets/cyber-threats-2019-retrospect.pdf\n  - https://www.fireeye.com/blog/threat-research/2018/09/apt10-targeting-japanese-corporations-using-updated-ttps.html\nauthor: SOC Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.defense_evasion\n  - attack.persistence\n  - attack.t1574.002\n  - attack.t1036.005\n  - detection.emerging_threats\nlogsource:\n  category: image_load\n  product: windows\ndetection:\n  selection_loader_name:\n    Image|endswith:\n      - '\\jjs.exe'\n      - '\\javaw.exe'\n      - '\\java.exe'\n      - '\\GoogleUpdate.exe'\n      - '\\mcoemcpy.exe'\n      - '\\mcods.exe'\n      - '\\DPMonitor.exe'\n      - '\\FlashHelp_Player.exe'\n      - '\\PotPlayerMini.exe'\n      - '\\vmnat.exe'\n      - '\\hkcmd.exe'\n  selection_loader_path:\n    Image|contains:\n      - '\\ProgramData\\'\n      - '\\AppData\\Local\\'\n      - '\\AppData\\Roaming\\'\n      - '\\Users\\Public\\'\n      - '\\PerfLogs\\'\n      - '\\Windows\\Temp\\'\n      - '\\Intel\\'\n  selection_dll:\n    ImageLoaded|endswith:\n      - '\\jli.dll'\n      - '\\msvcr100.dll'\n      - '\\mscoree.dll'\n      - '\\libcurl.dll'\n      - '\\oci.dll'\n      - '\\secur32.dll'\n      - '\\cryptbase.dll'\n      - '\\version.dll'\n      - '\\wtsapi32.dll'\n      - '\\netapi32.dll'\n      - '\\goopdate.dll'\n      - '\\FlashHelp.dll'\n      - '\\PotPlayer.dll'\n      - '\\vmtools.dll'\n  filter_system_paths:\n    Image|startswith:\n      - 'C:\\Windows\\'\n      - 'C:\\Program Files\\'\n      - 'C:\\Program Files (x86)\\'\n  condition: all of selection_* and not filter_system_paths\nfields:\n  - Image\n  - ImageLoaded\n  - User\n  - Hashes\n  - Signed\nfalsepositives:\n  - Portable application installers running out of staging directories\n  - Java JDK installations placed in non-default user directories\nlevel: high",
      "rationale": "Generic side-load detections fire on portable apps. This is locked to APT10 by intersecting THREE actor-specific facets: (1) the exact set of legit binaries Stone Panda has been observed re-purposing across the Cloud Hopper, A41APT and Operation Cuckoobees clusters (Java jjs/javaw for ANEL, mcoemcpy for ChChes/RedLeaves, FlashHelp_Player for UPPERCUT, GoogleUpdate.exe for PlugX), (2) APT10's docum"
    },
    {
      "title": "[LLM] APT10 / Cloudhopper trusted-relationship pivot \u2014 RDP from MSP / external host followed by schtasks persistence",
      "description": "Hunts the APT10 'Cloud Hopper' lateral-movement signature: an interactive RDP logon arriving on an internal host (often originating from an MSP-owned subnet, a jumpbox or a non-corporate IP) followed within 30 minutes by the same account spawning schtasks.exe /create to register a persistence task pointing at a loader in a writable directory (ProgramData / AppData / Public) \u2014 frequently masquerading as a Google/Windows updater. This RDP-then-schtasks chain is documented in PwC, BAE Systems and US-CERT TA17-117A on Cloud Hopper / Operation Cloud Hopper.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1199",
        "T1021.001",
        "T1053.005",
        "T1078",
        "T1036.005"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats `security_content_summariesonly` count min(_time) as logon_time from datamodel=Authentication where Authentication.action=success (Authentication.Logon_Type=10 OR Authentication.authentication_method=\"RemoteInteractive\" OR Authentication.signature_id=4624) by Authentication.dest Authentication.user Authentication.src _time | `drop_dm_object_name(Authentication)` | where NOT match(user, \"\\$$\") | rename _time as logon_time | join type=inner dest user [ | tstats `security_content_summariesonly` count from datamodel=Endpoint.Processes where Processes.process_name=\"schtasks.exe\" (Processes.process=\"*/create*\" OR Processes.process=\"*-create*\") by Processes.dest Processes.user Processes.process Processes.parent_process_name Processes.process_id _time | `drop_dm_object_name(Processes)` | rename _time as task_time ] | eval delay_sec = task_time - logon_time | where delay_sec >= 0 AND delay_sec <= 1800 | rex field=process \"(?i)/tr\\s+[\\\"']?(?<task_target>[^\\\"']+)\" | where match(task_target, \"(?i)\\\\\\\\(ProgramData|AppData|Users\\\\\\\\Public|PerfLogs|Temp)\\\\\\\\\") OR match(task_target, \"(?i)(powershell|rundll32|regsvr32|wscript|cscript|mshta|certutil|bitsadmin)\") | rex field=process \"(?i)/tn\\s+[\\\"']?(?<task_name>[^\\\"']+)\" | where match(task_name, \"(?i)(google|update|chrome|adobe|windows|defender|onedrive|java)\") OR match(task_target, \"(?i)\\\\.(exe|dll|bat|vbs|ps1|js)$\") | table logon_time task_time delay_sec dest user src task_name task_target process parent_process_name | sort - task_time",
      "kql": "// APT10 / Cloud Hopper \u2014 RDP arrival followed by schtasks persistence pointing at a writable directory\nlet WindowMin = 30m;\nlet RdpArrivals =\n    DeviceLogonEvents\n    | where Timestamp > ago(7d)\n    | where ActionType == \"LogonSuccess\"\n    | where LogonType == \"RemoteInteractive\"        // Logon Type 10 \u2014 RDP\n    | where AccountName !endswith \"$\"\n    | where AccountName !in~ (\"system\",\"local service\",\"network service\")\n    | project LogonTime = Timestamp, DeviceId, DeviceName,\n              ActorAccount = tolower(AccountName),\n              SourceIP = RemoteIP, SourceHost = RemoteDeviceName,\n              SourceIPType = RemoteIPType, LogonId;\nDeviceProcessEvents\n| where Timestamp > ago(7d)\n| where FileName =~ \"schtasks.exe\"\n| where ProcessCommandLine has_any (\"/create\", \"-create\")\n| where AccountName !endswith \"$\"\n| join kind=inner RdpArrivals on DeviceId\n| where Timestamp between (LogonTime .. LogonTime + WindowMin)\n| where tolower(AccountName) == ActorAccount\n| extend TaskTarget = extract(@\"(?i)/tr\\s+[\"\"']?([^\"\"']+)\", 1, ProcessCommandLine),\n         TaskName   = extract(@\"(?i)/tn\\s+[\"\"']?([^\"\"']+)\", 1, ProcessCommandLine),\n         RunAs      = extract(@\"(?i)/ru\\s+[\"\"']?([^\"\"']+)\", 1, ProcessCommandLine)\n| where TaskTarget has_any (@\"\\ProgramData\\\", @\"\\AppData\\\", @\"\\Users\\Public\\\", @\"\\PerfLogs\\\", @\"\\Windows\\Temp\\\", @\"\\Intel\\\")\n      or TaskTarget has_any (\"powershell\",\"pwsh\",\"rundll32\",\"regsvr32\",\"wscript\",\"cscript\",\"mshta\",\"certutil\",\"bitsadmin\")\n| extend MasqueradeTaskName = TaskName matches regex @\"(?i)(google|update|chrome|adobe|windows|defender|onedrive|java|edge)\"\n| project LogonTime, TaskCreationTime = Timestamp,\n          DelaySec = datetime_diff('second', Timestamp, LogonTime),\n          DeviceName, AccountName,\n          SourceIP, SourceIPType, SourceHost,\n          TaskName, TaskTarget, RunAs, MasqueradeTaskName,\n          ProcessCommandLine, ParentProcess = InitiatingProcessFileName\n| order by TaskCreationTime desc",
      "sentinel_kql": "// APT10 / Cloud Hopper \u2014 RDP arrival followed by schtasks persistence to writable directory (Sentinel SecurityEvent)\nlet WindowMin = 30m;\nlet RdpArrivals =\n    SecurityEvent\n    | where TimeGenerated > ago(7d)\n    | where EventID == 4624 and LogonType == 10           // RDP success\n    | where TargetUserName !endswith \"$\"\n    | where TargetUserName !in~ (\"system\",\"local service\",\"network service\",\"anonymous logon\")\n    | project LogonTime = TimeGenerated, Computer,\n              ActorAccount = tolower(TargetUserName),\n              SourceIP = IpAddress, SourceWorkstation = WorkstationName,\n              TargetLogonId;\nSecurityEvent\n| where TimeGenerated > ago(7d)\n| where EventID == 4688\n| where NewProcessName endswith \"\\\\schtasks.exe\"\n| where CommandLine has_any (\"/create\",\"-create\")\n| where Account !endswith \"$\"\n| join kind=inner RdpArrivals on Computer\n| where TimeGenerated between (LogonTime .. LogonTime + WindowMin)\n| where tolower(TargetUserName) == ActorAccount\n| extend TaskTarget = extract(@\"(?i)/tr\\s+[\"\"']?([^\"\"']+)\", 1, CommandLine),\n         TaskName   = extract(@\"(?i)/tn\\s+[\"\"']?([^\"\"']+)\", 1, CommandLine),\n         RunAs      = extract(@\"(?i)/ru\\s+[\"\"']?([^\"\"']+)\", 1, CommandLine)\n| where TaskTarget has_any (@\"\\ProgramData\\\", @\"\\AppData\\\", @\"\\Users\\Public\\\", @\"\\PerfLogs\\\", @\"\\Windows\\Temp\\\", @\"\\Intel\\\")\n     or TaskTarget has_any (\"powershell\",\"pwsh\",\"rundll32\",\"regsvr32\",\"wscript\",\"cscript\",\"mshta\",\"certutil\",\"bitsadmin\")\n| extend MasqueradeTaskName = TaskName matches regex @\"(?i)(google|update|chrome|adobe|windows|defender|onedrive|java|edge)\"\n| project LogonTime, TaskCreationTime = TimeGenerated,\n          DelaySec = datetime_diff('second', TimeGenerated, LogonTime),\n          Computer, Account = TargetUserName,\n          SourceIP, SourceWorkstation,\n          TaskName, TaskTarget, RunAs, MasqueradeTaskName,\n          CommandLine, ParentProcess = ParentProcessName\n| order by TaskCreationTime desc",
      "sigma_yaml": "",
      "rationale": "Schtasks /create alone is noisy (legit installers, IT scripts). Pivoting purely on RDP is also noisy. APT10's Cloud Hopper signature is the SEQUENCE: trusted-relationship inbound RDP (often from an MSP jumpbox or partner subnet) \u2192 same-session schtasks /create \u2192 task target in a writable directory or invoking a script-host LOLBin \u2192 task name masquerading as Google/Adobe/Windows updater. PwC and BA"
    }
  ]
}