# intel/.llm_actor_uc_cache/08/088d3a0bec37d5225692b85a182c25617b69d86d.json
{
  "name": "LazyScripter",
  "ucs": [
    {
      "title": "[LLM] LazyScripter (G0140) \u2014 open-source toolkit retrieval from GitHub raw / paste-site staging via LOLBin script interpreters",
      "description": "Hunts the signature LazyScripter staging behaviour: a Windows script-interpreter LOLBin (mshta, powershell, wscript/cscript, bitsadmin, certutil, curl, rundll32, regsvr32) calling raw.githubusercontent.com / pastebin / paste.ee / anonfiles to pull the next stage of an open-source RAT bundle (Koadic, Octopus, LuminosityLink, RMS, KOADIC, Quasar, njRat). Malwarebytes (Feb 2021) and follow-up CrowdStrike/MSFT reporting both pin GitHub-raw-hosted toolkits as the through-line for this airline-targeting actor \u2014 almost no other named cluster runs droppers that hard-code raw.githubusercontent.com URLs",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "delivery",
      "conf": "high",
      "techs": [
        "T1102",
        "T1105",
        "T1583.006",
        "T1218.005",
        "T1059.001",
        "T1059.005",
        "T1027.010"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as firstTime max(_time) as lastTime values(Processes.process) as cmdline values(Processes.parent_process_name) as parent values(Processes.parent_process) as parent_cmdline from datamodel=Endpoint.Processes where Processes.process_name IN (\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"wscript.exe\",\"cscript.exe\",\"bitsadmin.exe\",\"curl.exe\",\"certutil.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"cmd.exe\") AND (Processes.process=\"*raw.githubusercontent.com*\" OR Processes.process=\"*gist.githubusercontent.com*\" OR Processes.process=\"*pastebin.com*\" OR Processes.process=\"*paste.ee*\" OR Processes.process=\"*anonfiles.com*\" OR Processes.process=\"*transfer.sh*\" OR Processes.process=\"*controlc.com*\" OR Processes.process=\"*cdn.discordapp.com*\") AND NOT (Processes.user IN (\"*$\",\"SYSTEM\",\"LOCAL SERVICE\",\"NETWORK SERVICE\")) by host Processes.user Processes.process_name Processes.dest Processes.process_id Processes.parent_process_id\n| `drop_dm_object_name(Processes)`\n| `security_content_ctime(firstTime)`\n| `security_content_ctime(lastTime)`\n| where count < 25",
      "kql": "let LookbackDays = 14d;\n// LazyScripter (G0140) \u2014 open-source RAT staging hosts\nlet _StagingHosts = dynamic([\n    \"raw.githubusercontent.com\",\"gist.githubusercontent.com\",\n    \"pastebin.com\",\"paste.ee\",\"anonfiles.com\",\n    \"transfer.sh\",\"controlc.com\",\"cdn.discordapp.com\"]);\nlet _LOLBins = dynamic([\n    \"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"wscript.exe\",\"cscript.exe\",\n    \"bitsadmin.exe\",\"curl.exe\",\"certutil.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"cmd.exe\"]);\n// (1) cmdline-embedded URL \u2014 the canonical `mshta http://raw.githubusercontent.com/...` shape\nlet ProcessHits =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where FileName in~ (_LOLBins)\n    | where ProcessCommandLine has_any (_StagingHosts)\n    | project Timestamp, DeviceName, AccountName,\n              Parent = InitiatingProcessFileName,\n              ParentCmd = InitiatingProcessCommandLine,\n              Child = FileName,\n              ChildCmd = ProcessCommandLine,\n              SHA256, Source = \"ProcessCmd\";\n// (2) actual outbound socket from the LOLBin to the staging host\nlet NetworkHits =\n    DeviceNetworkEvents\n    | where Timestamp > ago(LookbackDays)\n    | where InitiatingProcessFileName in~ (_LOLBins)\n    | where InitiatingProcessAccountName !endswith \"$\"\n    | where RemoteUrl has_any (_StagingHosts)\n    | where RemoteIPType == \"Public\"\n    | project Timestamp, DeviceName,\n              AccountName = InitiatingProcessAccountName,\n              Parent = InitiatingProcessParentFileName,\n              ParentCmd = \"\",\n              Child = InitiatingProcessFileName,\n              ChildCmd = strcat(InitiatingProcessCommandLine, \" -> \", RemoteUrl),\n              SHA256 = InitiatingProcessSHA256,\n              Source = \"NetEgress\";\nunion ProcessHits, NetworkHits\n| order by Timestamp desc",
      "sentinel_kql": "let LookbackDays = 14d;\nlet _StagingHosts = dynamic([\n    \"raw.githubusercontent.com\",\"gist.githubusercontent.com\",\n    \"pastebin.com\",\"paste.ee\",\"anonfiles.com\",\n    \"transfer.sh\",\"controlc.com\",\"cdn.discordapp.com\"]);\nlet _LOLBins = dynamic([\n    \"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"wscript.exe\",\"cscript.exe\",\n    \"bitsadmin.exe\",\"curl.exe\",\"certutil.exe\",\"rundll32.exe\",\"regsvr32.exe\",\"cmd.exe\"]);\nSecurityEvent\n| where TimeGenerated > ago(LookbackDays)\n| where EventID == 4688\n| where Account !endswith \"$\"\n| extend BinName = tolower(tostring(split(NewProcessName, \"\\\\\")[-1]))\n| where BinName in (_LOLBins)\n| where CommandLine has_any (_StagingHosts)\n| project TimeGenerated, Computer, Account, NewProcessName, CommandLine,\n          ParentProcessName, TokenElevationType\n| order by TimeGenerated desc",
      "sigma_yaml": "title: LazyScripter (G0140) LOLBin Retrieval From GitHub/Paste-Site Staging\nid: 7b3a91c4-1cab-44e0-9d4f-lazyscript01\nstatus: experimental\ndescription: |\n  Detects the signature LazyScripter (G0140 / Malwarebytes Feb-2021) staging\n  technique: Windows script-interpreter LOLBins fetching open-source RAT\n  toolkits (Koadic, Octopus, LuminosityLink, RMS, Quasar, njRat) from\n  raw.githubusercontent.com / gist / pastebin / paste.ee / anonfiles. The\n  combination of script LOLBin + GitHub-raw URL on cmdline is rare in legit\n  enterprise software but a hard-coded pattern in LazyScripter droppers.\nreferences:\n  - https://attack.mitre.org/groups/G0140/\n  - https://www.malwarebytes.com/blog/news/2021/02/lazyscripter-from-empire-to-double-rat\nauthor: Detection Engineering\ndate: 2026/05/03\ntags:\n  - attack.t1102\n  - attack.t1105\n  - attack.t1218.005\n  - attack.t1059.001\n  - attack.t1059.005\n  - attack.g0140\nlogsource:\n  category: process_creation\n  product: windows\ndetection:\n  selection_image:\n    Image|endswith:\n      - '\\mshta.exe'\n      - '\\powershell.exe'\n      - '\\pwsh.exe'\n      - '\\wscript.exe'\n      - '\\cscript.exe'\n      - '\\bitsadmin.exe'\n      - '\\curl.exe'\n      - '\\certutil.exe'\n      - '\\rundll32.exe'\n      - '\\regsvr32.exe'\n      - '\\cmd.exe'\n  selection_url:\n    CommandLine|contains:\n      - 'raw.githubusercontent.com'\n      - 'gist.githubusercontent.com'\n      - 'pastebin.com'\n      - 'paste.ee'\n      - 'anonfiles.com'\n      - 'transfer.sh'\n      - 'controlc.com'\n      - 'cdn.discordapp.com/attachments'\n  filter_machine_acct:\n    User|endswith: '$'\n  condition: selection_image and selection_url and not filter_machine_acct\nfalsepositives:\n  - Developer machines pulling boilerplate from public gists\n  - DevOps bootstrap scripts (rare on user endpoints; baseline + allow-list)\nlevel: high",
      "rationale": "LazyScripter is the rare named cluster whose tradecraft *is* GitHub-raw staging \u2014 every Malwarebytes/MSFT writeup of G0140 highlights raw.githubusercontent.com URLs hard-coded inside their .bat/.vbs/.lnk droppers calling mshta or powershell. Generic 'LOLBin downloads from internet' rules are noisy; pinning the destination to a small set of public-paste/code-hosting domains AND requiring the actor'",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] LazyScripter (G0140) \u2014 IATA/airline-themed lure script dropper chain to Run-key persistence (\u226415 min)",
      "description": "Correlates the second LazyScripter signature: an Office/archive-utility/explorer.exe parent spawning wscript/cscript/mshta/powershell against a file in an Outlook attachment / Downloads / %TEMP% path (.vbs/.js/.hta/.bat/.lnk/.wsf), followed within 15 minutes on the same host by a HKCU\\\u2026\\Run or RunOnce value being set by another script LOLBin. This is the canonical LazyScripter drop-and-pin chain \u2014 zipped IATA/visa/airline-job lures \u2192 embedded VBS/HTA \u2192 Run-key persistence for their Octopus/KOADIC stagers \u2014 and the parent + persistence pair is what makes it actor-specific rather than generic 'w",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "medium",
      "techs": [
        "T1566.001",
        "T1204.002",
        "T1059.005",
        "T1218.005",
        "T1547.001",
        "T1036"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=true count min(_time) as DropperTime values(Processes.process) as DropperCmd values(Processes.parent_process_name) as DropperParent values(Processes.process_name) as Dropper from datamodel=Endpoint.Processes where Processes.parent_process_name IN (\"explorer.exe\",\"outlook.exe\",\"winrar.exe\",\"7zfm.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\") AND Processes.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\") AND (Processes.process=\"*\\\\Outlook\\\\*\" OR Processes.process=\"*\\\\Content.Outlook\\\\*\" OR Processes.process=\"*\\\\Downloads\\\\*\" OR Processes.process=\"*\\\\AppData\\\\Local\\\\Temp\\\\*\" OR Processes.process=\"*.vbs*\" OR Processes.process=\"*.js*\" OR Processes.process=\"*.hta*\" OR Processes.process=\"*.bat*\" OR Processes.process=\"*.lnk*\" OR Processes.process=\"*.wsf*\") AND NOT Processes.user IN (\"*$\",\"SYSTEM\",\"LOCAL SERVICE\",\"NETWORK SERVICE\") by host Processes.user\n| `drop_dm_object_name(Processes)`\n| rename host as dest\n| join type=inner dest [\n    | tstats summariesonly=true min(_time) as PersistenceTime values(Registry.registry_path) as RegPath values(Registry.registry_value_name) as RegValueName values(Registry.registry_value_data) as RegValueData values(Registry.process_name) as RegWriter from datamodel=Endpoint.Registry where (Registry.registry_path=\"*\\\\CurrentVersion\\\\Run\\\\*\" OR Registry.registry_path=\"*\\\\CurrentVersion\\\\RunOnce\\\\*\") AND Registry.action=\"modified\" AND Registry.process_name IN (\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\") by host\n    | `drop_dm_object_name(Registry)`\n    | rename host as dest\n]\n| where PersistenceTime >= DropperTime AND PersistenceTime <= DropperTime + 900\n| eval DropperTime=strftime(DropperTime,\"%Y-%m-%dT%H:%M:%S\"), PersistenceTime=strftime(PersistenceTime,\"%Y-%m-%dT%H:%M:%S\")\n| table DropperTime, PersistenceTime, dest, user, DropperParent, Dropper, DropperCmd, RegWriter, RegPath, RegValueName, RegValueData",
      "kql": "let LookbackDays = 14d;\nlet _DropperBins  = dynamic([\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\"]);\nlet _LureParents  = dynamic([\"explorer.exe\",\"outlook.exe\",\"winrar.exe\",\"7zfm.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\"]);\nlet _LurePaths    = dynamic([@\"\\Outlook\\\", @\"\\Content.Outlook\\\", @\"\\Downloads\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Temp\\\"]);\nlet _ScriptExts   = @\"(?i)\\.(vbs|js|jse|hta|bat|cmd|lnk|wsf|wsh|ps1)\\b\";\n// (1) Lure-attachment script execution\nlet Droppers =\n    DeviceProcessEvents\n    | where Timestamp > ago(LookbackDays)\n    | where AccountName !endswith \"$\"\n    | where InitiatingProcessFileName in~ (_LureParents)\n    | where FileName in~ (_DropperBins)\n    | where ProcessCommandLine has_any (_LurePaths)\n         or ProcessCommandLine matches regex _ScriptExts\n    | project DropperTime = Timestamp, DeviceId, DeviceName, AccountName,\n              DropperParent = InitiatingProcessFileName,\n              Dropper       = FileName,\n              DropperCmd    = ProcessCommandLine,\n              DropperSHA256 = SHA256;\n// (2) Run / RunOnce key written by a script LOLBin within 15 minutes\nDeviceRegistryEvents\n| where Timestamp > ago(LookbackDays)\n| where ActionType in (\"RegistryValueSet\",\"RegistryKeyCreated\")\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\")\n| where InitiatingProcessFileName in~ (_DropperBins)\n| where InitiatingProcessAccountName !endswith \"$\"\n| join kind=inner Droppers on DeviceId\n| where Timestamp between (DropperTime .. DropperTime + 15m)\n| project DropperTime, PersistenceTime = Timestamp,\n          DeviceName, AccountName,\n          DropperParent, Dropper, DropperCmd, DropperSHA256,\n          PersistWriter = InitiatingProcessFileName,\n          PersistWriterCmd = InitiatingProcessCommandLine,\n          RegistryKey, RegistryValueName, RegistryValueData\n| order by DropperTime desc",
      "sentinel_kql": "let LookbackDays = 14d;\nlet _DropperBins  = dynamic([\"wscript.exe\",\"cscript.exe\",\"mshta.exe\",\"powershell.exe\",\"pwsh.exe\",\"cmd.exe\"]);\nlet _LureParents  = dynamic([\"explorer.exe\",\"outlook.exe\",\"winrar.exe\",\"7zfm.exe\",\"winword.exe\",\"excel.exe\",\"powerpnt.exe\"]);\nlet _LurePaths    = dynamic([@\"\\Outlook\\\", @\"\\Content.Outlook\\\", @\"\\Downloads\\\", @\"\\AppData\\Local\\Temp\\\", @\"\\Temp\\\"]);\nlet _ScriptExts   = @\"(?i)\\.(vbs|js|jse|hta|bat|cmd|lnk|wsf|wsh|ps1)\\b\";\nlet Droppers =\n    ImProcessCreate\n    | where TimeGenerated > ago(LookbackDays)\n    | where ActorUsername !endswith \"$\"\n    | where ParentProcessName has_any (_LureParents)\n    | extend ChildBin = tolower(tostring(split(TargetProcessName, \"\\\\\")[-1]))\n    | where ChildBin in (_DropperBins)\n    | where TargetProcessCommandLine has_any (_LurePaths)\n         or TargetProcessCommandLine matches regex _ScriptExts\n    | project DropperTime = TimeGenerated, DvcHostname, ActorUsername,\n              DropperParent = ParentProcessName,\n              Dropper       = ChildBin,\n              DropperCmd    = TargetProcessCommandLine;\nImRegistryEvent\n| where TimeGenerated > ago(LookbackDays)\n| where RegistryKey has_any (@\"\\CurrentVersion\\Run\", @\"\\CurrentVersion\\RunOnce\")\n| extend Writer = tolower(tostring(split(ActingProcessName, \"\\\\\")[-1]))\n| where Writer in (_DropperBins)\n| join kind=inner Droppers on $left.DvcHostname == $right.DvcHostname\n| where TimeGenerated between (DropperTime .. DropperTime + 15m)\n| project DropperTime, PersistenceTime = TimeGenerated,\n          DvcHostname, ActorUsername,\n          DropperParent, Dropper, DropperCmd,\n          PersistWriter = Writer,\n          RegistryKey, RegistryValue, RegistryValueData\n| order by DropperTime desc",
      "sigma_yaml": "",
      "rationale": "LazyScripter's airline-targeting tradecraft is unusually consistent: IATA/visa/job-application lure ZIPs containing a .vbs/.lnk/.hta dropper which Office/explorer/winrar opens, and the dropper drops a Run-key for the open-source RAT (KOADIC/Octopus/RMS) within seconds. Generic 'script from temp' or generic 'Run-key write' rules each generate hundreds of FPs/day; chaining the two with a 15-minute w"
    }
  ]
}