# intel/.llm_actor_uc_cache/3e/3e26dc8c02aafebcffb2878059739b5748196f0c.json
{
  "name": "UNC3886",
  "ucs": [
    {
      "title": "[LLM] UNC3886 Linux rootkit (REPTILE / MEDUSA) install via /etc/ld.so.preload hook",
      "description": "Hunts for writes to /etc/ld.so.preload (or files dropped under /etc/ld.so.preload.d/) by anything other than the OS package manager. Mandiant has tied this exact LD_PRELOAD persistence primitive to UNC3886's REPTILE and MEDUSA userland rootkits, deployed onto Linux guest VMs after the actor pivots from compromised vCenter / Fortinet appliances. Combined with a non-package-manager parent, this is one of the highest-fidelity signals of UNC3886 post-exploitation on a Linux estate.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1014",
        "T1554",
        "T1574",
        "T1070.004"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Filesystem.file_hash) as file_hash values(Filesystem.process_id) as proc_id from datamodel=Endpoint.Filesystem where (Filesystem.file_path=\"/etc/ld.so.preload\" OR Filesystem.file_path=\"/etc/ld.so.preload.d/*\") AND Filesystem.action IN (\"created\",\"modified\",\"write\",\"renamed\") AND NOT Filesystem.process_name IN (\"dpkg\",\"apt\",\"apt-get\",\"yum\",\"dnf\",\"rpm\",\"ldconfig\",\"zypper\",\"pacman\") by host Filesystem.dest Filesystem.file_path Filesystem.file_name Filesystem.process_name Filesystem.process Filesystem.user | `drop_dm_object_name(Filesystem)` | join type=outer host process_name [| tstats summariesonly=t values(Processes.parent_process_name) as parent_process_name values(Processes.process) as process_full from datamodel=Endpoint.Processes by host Processes.process_name | `drop_dm_object_name(Processes)`] | convert ctime(firstTime) ctime(lastTime) | table firstTime lastTime host dest file_path file_name file_hash process_name process parent_process_name user",
      "kql": "// UNC3886 (Mandiant) \u2014 REPTILE / MEDUSA rootkit persistence via dynamic linker preload hook\nlet _allowed_writers = dynamic([\"dpkg\",\"apt\",\"apt-get\",\"yum\",\"dnf\",\"rpm\",\"ldconfig\",\"zypper\",\"pacman\",\"unattended-upgrade\"]);\nlet _system_paths   = dynamic([@\"/usr/bin/\",@\"/usr/sbin/\",@\"/bin/\",@\"/sbin/\"]);\nDeviceFileEvents\n| where Timestamp > ago(7d)\n| where ActionType in (\"FileCreated\",\"FileModified\",\"FileRenamed\")\n| where FolderPath startswith \"/etc/\"                       // narrow before substring match\n| where (FolderPath == \"/etc/\" and FileName =~ \"ld.so.preload\")\n     or FolderPath startswith \"/etc/ld.so.preload.d/\"\n| where InitiatingProcessFileName !in~ (_allowed_writers)\n| extend NotFromPkgMgrPath = not(InitiatingProcessFolderPath has_any (_system_paths))\n| project Timestamp, DeviceName, ActionType, FolderPath, FileName,\n          SHA256,\n          Writer  = InitiatingProcessFileName,\n          WriterCmd = InitiatingProcessCommandLine,\n          WriterPath = InitiatingProcessFolderPath,\n          WriterAcct = InitiatingProcessAccountName,\n          WriterParent = InitiatingProcessParentFileName,\n          NotFromPkgMgrPath\n| order by Timestamp desc",
      "sentinel_kql": "// UNC3886 (Mandiant) \u2014 REPTILE / MEDUSA rootkit persistence via dynamic linker preload hook\nlet _allowed_writers = dynamic([\"dpkg\",\"apt\",\"apt-get\",\"yum\",\"dnf\",\"rpm\",\"ldconfig\",\"zypper\",\"pacman\",\"unattended-upgrade\"]);\nImFileEvent\n| where TimeGenerated > ago(7d)\n| where EventType in (\"FileCreated\",\"FileModified\",\"FileRenamed\")\n| where TargetFilePath startswith \"/etc/\"\n| where (TargetFilePath == \"/etc/ld.so.preload\")\n     or TargetFilePath startswith \"/etc/ld.so.preload.d/\"\n| where ActingProcessName !in~ (_allowed_writers)\n| project TimeGenerated, DvcHostname, EventType, TargetFilePath, TargetFileName,\n          TargetFileSHA256,\n          Writer = ActingProcessName,\n          WriterCmd = ActingProcessCommandLine,\n          WriterUser = ActorUsername\n| order by TimeGenerated desc",
      "sigma_yaml": "title: UNC3886 LD_PRELOAD Rootkit Persistence (REPTILE/MEDUSA)\nid: 6f3d3e1c-5b8f-4f0a-9a45-unc3886-ldpreload\nstatus: experimental\ndescription: Detects writes to /etc/ld.so.preload or files placed under /etc/ld.so.preload.d/ by a process that is not a known package manager. Tied to UNC3886 deployment of the REPTILE and MEDUSA Linux rootkits onto guest VMs reached from compromised vCenter / Fortinet appliances.\nreferences:\n  - https://cloud.google.com/blog/topics/threat-intelligence/uncovering-unc3886-espionage-operations\n  - https://attack.mitre.org/groups/G1048/\nauthor: SOC Detection Engineering\ntags:\n  - attack.persistence\n  - attack.defense_evasion\n  - attack.t1014\n  - attack.t1554\n  - attack.t1574\n  - detection.threat_hunting\nlogsource:\n  product: linux\n  category: file_event\ndetection:\n  selection_path:\n    TargetFilename|startswith:\n      - '/etc/ld.so.preload.d/'\n    TargetFilename|contains:\n      - '/etc/ld.so.preload'\n  filter_pkgmgr:\n    Image|endswith:\n      - '/dpkg'\n      - '/apt'\n      - '/apt-get'\n      - '/yum'\n      - '/dnf'\n      - '/rpm'\n      - '/ldconfig'\n      - '/zypper'\n      - '/pacman'\n      - '/unattended-upgrade'\n  condition: selection_path and not filter_pkgmgr\nfields:\n  - Image\n  - CommandLine\n  - TargetFilename\n  - User\n  - ParentImage\nfalsepositives:\n  - Custom in-house performance/security shims deployed via configuration management (Ansible/Salt/Puppet) \u2014 should be allowlisted by writer path or hash.\nlevel: high",
      "rationale": "REPTILE and MEDUSA are the two Linux userland rootkits Mandiant attributes specifically to UNC3886 after vCenter/ESXi compromise. Both rely on the LD_PRELOAD primitive \u2014 either dropping a malicious .so referenced from /etc/ld.so.preload or chaining one in /etc/ld.so.preload.d/. Writes to those paths are vanishingly rare on a healthy Linux estate, and the package-manager exclusion strips the only r",
      "_sigma_issues": [
        "sigma parse error: Sigma rule identifier must be an UUID"
      ]
    },
    {
      "title": "[LLM] UNC3886 ESXi/Linux rc.local persistence implant after SSH lateral movement from vCenter",
      "description": "Hunts modifications to /etc/rc.local, /etc/rc.local.d/local.sh, /etc/init.d/* and systemd unit files where the writer is a shell child of sshd (i.e. an interactive SSH session installing persistence). Mandiant documented UNC3886 backdooring exactly /etc/rc.local.d/local.sh on ESXi hosts to launch VIRTUALPITA / VIRTUALGATE at boot, after pivoting via SSH from a compromised vCenter. The sshd\u2192bash\u2192write-to-rc-local chain is the actor's signature install step.",
      "is_llm": true,
      "source_kind": "actor-bespoke",
      "phase": "install",
      "conf": "high",
      "techs": [
        "T1037.004",
        "T1021.004",
        "T1059.004",
        "T1505.006"
      ],
      "art_id": "",
      "art_title": "",
      "is_mitre_match": false,
      "splunk": "| tstats summariesonly=t count min(_time) as firstTime max(_time) as lastTime values(Filesystem.file_hash) as file_hash from datamodel=Endpoint.Filesystem where (Filesystem.file_path=\"/etc/rc.local\" OR Filesystem.file_path=\"/etc/rc.local.d/*\" OR Filesystem.file_path=\"/etc/init.d/*\" OR Filesystem.file_path=\"/etc/systemd/system/*\" OR Filesystem.file_path=\"/usr/local/etc/rc.d/*\") AND Filesystem.action IN (\"created\",\"modified\",\"write\",\"renamed\") AND Filesystem.process_name IN (\"bash\",\"sh\",\"dash\",\"ash\",\"zsh\",\"python\",\"python3\",\"perl\",\"vi\",\"vim\",\"nano\",\"sed\",\"tee\",\"cp\",\"mv\",\"scp\",\"sftp-server\") by host Filesystem.dest Filesystem.file_path Filesystem.file_name Filesystem.process_name Filesystem.process Filesystem.process_id Filesystem.user | `drop_dm_object_name(Filesystem)` | rename process_id as fs_pid | join type=inner host fs_pid [| tstats summariesonly=t values(Processes.parent_process_name) as parent_process_name values(Processes.parent_process) as parent_process_cmd values(Processes.process) as full_cmd from datamodel=Endpoint.Processes where Processes.parent_process_name=\"sshd\" by host Processes.process_id | `drop_dm_object_name(Processes)` | rename process_id as fs_pid] | convert ctime(firstTime) ctime(lastTime) | table firstTime lastTime host dest file_path file_name file_hash process_name full_cmd parent_process_name user",
      "kql": "// UNC3886 (Mandiant) \u2014 VIRTUALPITA / VIRTUALGATE rc.local persistence after SSH pivot from vCenter\nlet _persist_paths = dynamic([\n  \"/etc/rc.local\",\n  \"/etc/rc.local.d/\",\n  \"/etc/init.d/\",\n  \"/etc/systemd/system/\",\n  \"/usr/local/etc/rc.d/\"\n]);\nlet _shells = dynamic([\"bash\",\"sh\",\"dash\",\"ash\",\"zsh\",\"python\",\"python3\",\"perl\",\"vi\",\"vim\",\"nano\",\"sed\",\"tee\",\"cp\",\"mv\",\"scp\",\"sftp-server\"]);\n// Writers that are children of sshd within the past 7d\nlet _ssh_session_pids = DeviceProcessEvents\n    | where Timestamp > ago(7d)\n    | where InitiatingProcessFileName =~ \"sshd\"\n    | where FileName in~ (_shells)\n    | project DeviceId, ShellPid = ProcessId, ShellTime = Timestamp,\n              ShellAcct = AccountName, RemoteParent = InitiatingProcessParentFileName,\n              ShellCmd = ProcessCommandLine;\nDeviceFileEvents\n| where Timestamp > ago(7d)\n| where ActionType in (\"FileCreated\",\"FileModified\",\"FileRenamed\")\n| where FolderPath has_any (_persist_paths) or FileName =~ \"local.sh\"\n| where InitiatingProcessFileName in~ (_shells)\n| join kind=inner _ssh_session_pids on $left.DeviceId == $right.DeviceId\n        and $left.InitiatingProcessId == $right.ShellPid\n| where Timestamp between (ShellTime .. ShellTime + 1h)\n| project Timestamp, DeviceName, FolderPath, FileName, ActionType, SHA256,\n          Writer = InitiatingProcessFileName,\n          WriterCmd = InitiatingProcessCommandLine,\n          WriterAcct = InitiatingProcessAccountName,\n          ShellTime, ShellCmd, ShellAcct\n| order by Timestamp desc",
      "sentinel_kql": "// UNC3886 (Mandiant) \u2014 VIRTUALPITA / VIRTUALGATE rc.local persistence after SSH pivot from vCenter\nlet _persist_paths = dynamic([\n  \"/etc/rc.local\",\n  \"/etc/rc.local.d/\",\n  \"/etc/init.d/\",\n  \"/etc/systemd/system/\",\n  \"/usr/local/etc/rc.d/\"\n]);\nlet _shells = dynamic([\"bash\",\"sh\",\"dash\",\"ash\",\"zsh\",\"python\",\"python3\",\"perl\",\"vi\",\"vim\",\"nano\",\"sed\",\"tee\",\"cp\",\"mv\",\"scp\",\"sftp-server\"]);\nlet _ssh_shell_pids = ImProcessCreate\n    | where TimeGenerated > ago(7d)\n    | where ParentProcessName =~ \"sshd\"\n    | where ActingProcessName in~ (_shells)\n    | project DvcHostname, ShellPid = ActingProcessId, ShellTime = TimeGenerated,\n              ShellUser = ActorUsername, ShellCmd = ActingProcessCommandLine;\nImFileEvent\n| where TimeGenerated > ago(7d)\n| where EventType in (\"FileCreated\",\"FileModified\",\"FileRenamed\")\n| where TargetFilePath has_any (_persist_paths) or TargetFileName =~ \"local.sh\"\n| where ActingProcessName in~ (_shells)\n| join kind=inner _ssh_shell_pids on $left.DvcHostname == $right.DvcHostname\n        and $left.ActingProcessId == $right.ShellPid\n| where TimeGenerated between (ShellTime .. ShellTime + 1h)\n| project TimeGenerated, DvcHostname, TargetFilePath, TargetFileName, EventType, TargetFileSHA256,\n          Writer = ActingProcessName, WriterCmd = ActingProcessCommandLine,\n          ShellTime, ShellCmd, ShellUser\n| order by TimeGenerated desc",
      "sigma_yaml": "",
      "rationale": "UNC3886 specifically \u2014 not random APTs \u2014 chains SSH-from-vCenter into ESXi/Linux rc.local backdooring to autostart VIRTUALPITA and VIRTUALGATE on hypervisor boot. Mandiant's named writeups call out /etc/rc.local.d/local.sh as the implant location. Plenty of detections look at rc.local writes; few require the writer to be a sshd child within the same hour. That sshd\u2192shell\u2192persistence-file chain is "
    }
  ]
}